# 0-raw-data/ORKL/library/jsons_entries/0068eda0-1115-4877-aa72-81da839c451e.json
{
    "id": "0068eda0-1115-4877-aa72-81da839c451e",
    "created_at": "2023-01-12T15:05:34.321032Z",
    "updated_at": "2025-03-27T02:05:26.706634Z",
    "deleted_at": null,
    "sha1_hash": "d0af47724620456f4a09067901e453b0dd6d3b58",
    "title": "2021-09-21 - Cryptominer z0Miner Uses Newly Discovered Vulnerability CVE-2021-26084 to Its Advantage",
    "authors": "",
    "file_creation_date": "2022-05-28T15:26:40Z",
    "file_modification_date": "2022-05-28T15:26:40Z",
    "file_size": 324767,
    "plain_text": "# Cryptominer z0Miner Uses Newly Discovered Vulnerability CVE-2021-26084 to Its Advantage\n\n**[trendmicro.com/en_us/research/21/i/cryptominer-z0miner-uses-newly-discovered-vulnerability-cve-2021.html](https://www.trendmicro.com/en_us/research/21/i/cryptominer-z0miner-uses-newly-discovered-vulnerability-cve-2021.html)**\n\nSeptember 21, 2021\n\nRecently, we discovered that the cryptomining trojan z0Miner has been taking advantage of the Atlassian’s Confluence remote code\n[execution (RCE) vulnerability assigned as CVE-2021-26084, which was disclosed by Atlassian in August. Given the increasing](https://nvd.nist.gov/vuln/detail/CVE-2021-26084)\npopularity of the cryptocurrency market, we expect malware authors behind trojans like z0Miner to constantly update the techniques\nand entry vectors they use to gain a foothold within a system.\n\n[This trojan was initially observed exploiting Oracle’s WebLogic Server RCE, CVE-2020-14882, late last year. Since then, z0Miner has](https://nvd.nist.gov/vuln/detail/CVE-2020-14882)\nbeen gaining attention by utilizing different unauthorized RCE vulnerabilities, such as the ElasticSearch RCE bug, aka CVE-20151427.\n\n## Infection chain\n\nBased on our investigation, we found that the infection chain that leverages the new CVE-2021-26084 flaw (Figure 1) is identical to\nprevious findings on z0Miner, as reported by [360 Netlab and](https://blog.netlab.360.com/threat-alert-z0miner-is-spreading-quickly-by-exploiting-elasticsearch-and-jenkins-vulnerabilities/) [Tencent Security.](https://s.tencent.com/research/report/1170.html)\n\nOnce the Confluence vulnerability is successfully exploited, z0Miner deploys web shells that will download the following malicious files:\n\nhxxp://213[.]152[.]165[.]29/x[.]bat: detected by Trend Micro as Trojan.BAT.TINYOMED.ZYII\nhxxp://213[.]152[.]165[.]29/uninstall[.]bat: detected by Trend Micro as Trojan.BAT.SVCLAUNCHER.ZYII\nhxxp://213[.]152[.]165[.]29/vmicguestvs[.]dll: detected by Trend Micro as Trojan.Win64.TINYOMED.ZYII\nhxxp://27[.]1[.]1[.]34:8080/docs/s/sys[.]ps1: detected by Trend Micro as Trojan.PS1.Z0MINER.YXAIJ\n\nFigure 1. The infection chain of z0Miner\nEvasion mechanisms\n\n\n-----\n\nThe malware has been known to use several persistence and defense evasion mechanisms, one of which is the installation of the file\n_vmicvguestvs.dll that z0Miner disguises as a legitimate integration service called “Hyper-V Guest Integration” (Figure 2)._\n\nFigure 2. The creation\n\nof the fraudulent “Hyper-V Guest Integration” service\nOne of the downloaded scripts will also create a scheduled task called .NET Framework NGEN v4.0.30319 32 that poses as a .NET\nFramework NGEN task, as shown in Figure 3. This scheduled task is designed to download and execute a script from Pastebin every\nfive minutes. However, as of this writing, the contents of the Pastebin URL have already been taken down.\n\nFigure 3. The creation\n\nof the scheduled task\nThe z0Miner trojan will proceed to collect its own mining tools from URLs contained in the file ok.bat, as shown in Figure 4. It also\ndownloads another script named clean.bat to find and delete any cryptocurrency mining payloads from other competitors (Figure 5).\n\nFigure 4. The URLs\n\nand file paths of z0Miner’s mining components from the file ok.bat\n\nFigure 5. The clean.bat file that locates and deletes other cryptominers\n\nSecurity recommendations\n\nAlthough Atlassian has already released a patch addressing the Confluence vulnerability, users can take further steps to minimize their\nsystem’s exposure to threats like z0Miner. Regularly updating their systems and applications with the latest patches plays a critical role\nin mitigating the risks for end-users, ensuring that these security gaps can't be abused for malicious activities.\n\n[To assist with patch management, users can turn to solutions such as Trend Micro™ Deep Security™ and](https://www.trendmicro.com/en_us/business/products/hybrid-cloud/deep-security.html) Trend Micro Cloud One™ Workload Security, which provide virtual patching that protects servers and endpoints from threats that abuse vulnerabilities in critical\n[applications. Trend Micro ™ Deep Discovery™ offers detection, in-depth analysis, and a proactive response to attacks using exploits](https://www.trendmicro.com/en_us/business/products/network/advanced-threat-protection.html)\nand other similar threats through specialized engines, custom sandboxing, and seamless correlation across the entire attack life cycle,\nallowing it to detect threats even without any engine or pattern update.\n\nSimilarly, Workload Security defends systems and detects vulnerabilities and malware with the broadest hybrid cloud security\ncapabilities for a mixed environment of virtual, physical, cloud, and containers. Using techniques like machine learning (ML) and virtual\npatching, Workload Security also protects new and existing workloads even against unknown threats. It also shields users from\nexploits that target the Confluence vulnerability via the following rule:\n\n1011117 - Atlassian Confluence Server Remote Code Execution Vulnerability (CVE-2021-26084)\n\n[Users can also benefit from theTippingPoint® Threat Protection System, which uses comprehensive and contextual awareness](https://www.trendmicro.com/en_us/business/products/network/intrusion-prevention/tipping-point-threat-protection-system.html)\nanalysis for advanced threats that exploit vulnerabilities. Threat intelligence from sources such as Digital Vaccine Labs (DVLabs) and\nZero Day Initiative (ZDI) provides maximum threat coverage and virtual patching shields vulnerabilities against exploits. TippingPoint\nprotects customers through the following rule:\n\n40260: HTTP: Atlassian Confluence Server and Data Center OGNL Injection Vulnerability\n\n\n-----\n\nMITRE ATT&CK Tactics and Techniques\n\nThe following are the MITRE ATT&CK tactics and techniques associated with CVE-2021-26084 bundled with z0Miner:\n\n**Tactic** **Technique**\n\nExecution T1569.002: System Services: Service Execution\n\nPersistence T1053.005: Scheduled Task\n\nT1543.003: Create or Modify System Process: Systems Service\n\nDefense Evasion T1112: Modify Registry\n\nT1489: Service Stop\n\nT1562.001: Impair Defenses: Disable or Modify Tools\n\nT1036.004: Masquerade Task or Service\n\nT1070.004: File Deletion\n\nDiscovery T1033: System Owner/User Discovery\n\nT1049: System Network Connections Discovery\n\nT1069.001: Permission Groups Discovery: Local Groups\n\nT1069.002: Permission Groups Discovery: Domain Groups\n\nT1082: System Information Discovery\n\nT1087: Account Discovery\n\nT1087.001: Account Discovery: Local Account\n\nT1087.002: Account Discovery: Domain Account\n\nT1124: System Time Discovery\n\nImpact T1496: Resource Hijacking\n\nIndicators of compromise\n\n**SHA-256** **Filename** **Trend Micro Detection Name**\n\n49f3d06419d9578551e584515f44b2ee714e1eef96b94e68ea957f2943deca5a error.jsp Possible_SMASPWEBSHELL\n\ncb339d08c0ad7c4d07b06cae5d7eae032fb1bb1178d80b2a1997a8b8257b5bea uninstall.bat Backdoor.Java.WEBSHELL.SBJKTK\n\n0663d70411a20340f184ae3b47138b33ac398c800920e4d976ae609b60522b01 _wxm.exe_ PUA.Win64.Xmrig.KBL\n\na5604893608cf08b7cbfb92d1cac20868808218b3cc453ca86da0abaeadc0537 _network02.exe_ Coinminer.Win64.MALXMR.SMA\n\nf176d69f18cde008f1998841c343c3e5d4337b495132232507a712902a0aec5e _.solrg_ Trojan.SH.Z0MINER.YXAIJ\n\n4a2fbe904e4665939d8517c48fb3d5cb67e9b1482195c41fe31396318118cfc8 _sys.ps1_ Trojan.PS1.Z0MINER.YXAIJ\n\ne9ba929949c7ea764a298e33af1107ff6feefe884cabf6254ff574efff8a2e40\n\n7d8b52e263bc548891c1623695bac7fb21dab112e43fffb515447a5cc709ac89 _clean.bat_ Trojan.BAT.KILLMINE.YXAIJ\n\n**URLs**\n\nhxxp://209.141.40.190/oracleservice.exe\nhxxp://209.141.40.190/wxm.exe\nhxxp://27.1.1.34:8080/docs/s/config.json\nhxxp://27.1.1.34:8080/examples/clean.bat\nhxxp://27.1.1.34:8080/docs/s/sys.ps1\nhxxp://222.122.47.27:2143/auth/xmrig.exe\nhxxp://pastebin.com/raw/bcFqDdXx\nhxxp://pastebin com/raw/g93wWHkR\n\n\n-----\n\nhxxp://164.52.212.196:88/eth.jpg\nhxxp://66.42.117.168/BootCore_jsp\nhxxp://164.52.212.196:88/1.jpg\nhxxp://209.141.40.190/xms\nhxxp://172.96.249.219:88/.jpg\nhxxp://172.96.249.219:88/1.jpg 1.bat\nhxxp://172.96.249.219:88/.jpg\nhxxps://zgpay.cc/css/kwork.sh\nhxxps://raw.githubusercontent.com/alreadyhave/thinkabout/main/kwork.sh\nhxxp://209.141.40.190/oracleservice.exe\nhxxp://213.152.165.29/vmicguestvs.dll\nhxxp://213.152.165.29/uninstall.bat\nhxxp://213.152.165.29/x.bat\n\nExploits & Vulnerabilities\n\nRecently, we discovered that the cryptomining trojan z0Miner has been taking advantage of the Atlassian’s Confluence remote code\nexecution (RCE) vulnerability assigned as CVE-2021-26084, which was disclosed by Atlassian in August.\n\nBy: Nikki Madayag, Josefino Fajilago IV September 21, 2021 Read time: ( words)\n\nContent added to Folio\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2021/2021-09-21 - Cryptominer z0Miner Uses Newly Discovered Vulnerability CVE-2021-26084 to Its Advantage.pdf"
    ],
    "report_names": [
        "2021-09-21 - Cryptominer z0Miner Uses Newly Discovered Vulnerability CVE-2021-26084 to Its Advantage.pdf"
    ],
    "threat_actors": [],
    "ts_created_at": 1673535934,
    "ts_updated_at": 1743041126,
    "ts_creation_date": 1653751600,
    "ts_modification_date": 1653751600,
    "files": {
        "pdf": "https://archive.orkl.eu/d0af47724620456f4a09067901e453b0dd6d3b58.pdf",
        "text": "https://archive.orkl.eu/d0af47724620456f4a09067901e453b0dd6d3b58.txt",
        "img": "https://archive.orkl.eu/d0af47724620456f4a09067901e453b0dd6d3b58.jpg"
    }
}