# 0-raw-data/ORKL/library/jsons_entries/0fc60245-3425-4adf-a69e-7ea4efc10f40.json
{
    "id": "0fc60245-3425-4adf-a69e-7ea4efc10f40",
    "created_at": "2023-01-12T15:06:40.828071Z",
    "updated_at": "2025-03-27T02:05:39.218607Z",
    "deleted_at": null,
    "sha1_hash": "d75cf01f6361a7f6e109e408656a8940b1a21c95",
    "title": "2021-06-03 - WebLogic RCE Leads to XMRig",
    "authors": "",
    "file_creation_date": "2022-05-28T18:03:53Z",
    "file_modification_date": "2022-05-28T18:03:53Z",
    "file_size": 1450485,
    "plain_text": "# WebLogic RCE Leads to XMRig\n\n**[thedfirreport.com/2021/06/03/weblogic-rce-leads-to-xmrig/](https://thedfirreport.com/2021/06/03/weblogic-rce-leads-to-xmrig/)**\n\n## Intro\n\n\nJune 3, 2021\n\n\nThis report will review an intrusion where, the threat actor took advantage of a WebLogic\nremote code execution vulnerability (CVE-2020–14882) to gain initial access to the system\nbefore installing a coin miner (XMRig).\n\n## Summary\n\nIn the attack, the threat actor took around 2 minutes from initial access, to running a\npersistent coin miner. The speed and number of times the exploit was ran against the target\nmake it likely that this was an automated attack, rather than human operated intrusion.\n\nThe threat actor exploited CVE-2020-14882 by making a request to the WebLogic server,\nwhich allowed the attacker to execute code. The payload then ran a PowerShell command to\ndownload a PowerShell script from a remote server. This script downloaded a loader and\nexecuted it, then it enabled persistence mechanisms. Next, the loader downloaded XMRig,\nand started the mining process.\n\nWhile this attack took place around a month ago the payloads used in the intrusion continue\nto be hosted as of this publication, indicating that this threat is still active and finding\nvulnerable WebLogic hosts to continue exploiting.\n\n\n-----\n\n[Analysis and reporting completed by @tas_kmanager and](https://twitter.com/tas_kmanager) [@TheDFIRReport](https://twitter.com/TheDFIRReport)\n\n[Reviewed by @iiamaleks and](https://twitter.com/iiamaleks) [@samaritan_o](https://twitter.com/samaritan_o)\n\n## MITRE ATT&CK\n\n Initial Access\n\nBefore jumping into the details of the attack, it is pertinent to understand CVE-2020-14882.\nThis vulnerability was discovered back in October 2020, which effects Oracle’s WebLogic\n[product. [1] We previously reported on an intrusion that also used this exploit and ended in a](https://www.oracle.com/security-alerts/cpuoct2020traditional.html)\ncoin miner. This vulnerability is easy to exploit, and leads to remote code execution (RCE)\nwithout authentication.\n\nJust by invoking URL path /console/images/%252E%252E%252Fconsole.portal and a little\nbit of tweaking of the packet sent to the server, the threat actor will be able to run their code\non the server.\n\nIn this case, we can clearly see from network logs that the attacker was probing the portal\npage and followed that by sending the crafted URL request.\n\nThe detailed payload in the header can be observed below:\n\nInterestingly, we can see the threat actor chose the FileSystemXmlApplicationContext class\nto carry out this attack, executing a file named poc.xml. This technique is listed in PoC\n[developed by Jang. [2]](https://github.com/jas502n/CVE-2020-14882#filesystemxmlapplicationcontext)\n\nBelow we can see the payload inside of the xml file hosted in the attacker server. The\npayload will execute the PowerShell command on the vulnerable system.\n\n\n-----\n\n## Execution\n\nThe threat actor leveraged the WebLogic vulnerability to spawn a command shell from the\nserver running in the Java process, which then in turn, was used to run PowerShell and\ncollect the final payload that the threat actor wished to run on the system.\n\nNext, let us investigate the content of the PowerShell script (ldr.ps1)\n\nOn line one, we can see var $cc (can be pronounced as C2?) that indicates the server where\nit will get all the other components such as the payload sys.exe. The next line sets up\nrandom characters to be appended to the payload executable.\n\nThis binary (renamed to sysvr013.exe) was then executed by the same script.\n\n\n-----\n\n## Persistence\n\nThere are several persistence mechanisms observed. The first one utilizes schtasks to\ncreate a task called “BrowserUpdate”.\n```\nschtasks /create /F /sc minute /mo 1 /tn \"BrowserUpdate\" /tr \"$dst\"\n\n```\nIt also created a registry run key to execute the miner binary on reboot.\n```\nreg add HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run /v Run /d \"$dst\" /t REG_SZ\n/f\n\n```\nWe also observed the Loader binary adding npf.sys (WinPcap binary) as a scheduled task.\n\n## Privilege Escalation\n\nNo privilege escalation was observed, but note that the Java.exe process was running as a\n**high integrity process from the start.**\n\n## Defense Evasion\n\nFrom the PowerShell script above (ldr.ps1), we can see that name is being randomized to\navoid detection by blue team. We can also see that the attacker disabled firewall rules to\nmake sure connection to the mining pool was not blocked.\n```\nnetsh advfirewall set allprofiles state off\n\n```\nAs noted above, the attacker created a task called “BrowserUpdate” to masquerade the\nscheduled task as a web browser update task.\n\n## Discovery\n\nWe can see several local discovery commands executed from the PowerShell script. They\nare looking for existing mining process on the machine, such as kthreaddi, sysrv, etc. They\nalso utilize netstat to find if any process is using mining related ports, such as 3333, 4444,\netc.\n\n\n-----\n\nThis is likely partially to kill other rival miners but also to prevent duplicate miners on the\nsame host as the exploit was executed against the target repeatedly in quick succession, at\nleast 12 times.\n\n## Impact\n\nWe saw that the additional sys.exe binary and ldr.ps1 PowerShell script was brought in to\nsupport the Mining operation from IP address 194[.]145[.]227[.]21\n\nWe also observed connection from [kthreaddi].exe talking to crypto mining pools\n\n\n-----\n\nEventually, the crypto mining operation will cause performance issue with the WebLogic\nserver.\n\n## XMRig Config file\n\n\n-----\n\n```\n{\n\"api\": {\n\"id\": null,\n\"worker-id\": null\n},\n\"autosave\": false,\n\"background\": true,\n\"randomx\": {\n\"init\": -1,\n\"mode\": \"auto\",\n\"1gb-pages\": false,\n\"rdmsr\": true,\n\"wrmsr\": true,\n\"cache_qos\": false,\n\"numa\": true,\n\"scratchpad_prefetch_mode\": 1\n},\n\"cpu\": {\n\"enabled\": true,\n\"huge-pages\": true,\n\"huge-pages-jit\": false,\n\"hw-aes\": null,\n\"priority\": null,\n\"memory-pool\": false,\n\"yield\": true,\n\"max-threads-hint\": 100,\n\"asm\": true,\n\"argon2-impl\": null,\n\"astrobwt-max-size\": 550,\n\"cn/0\": false,\n\"cn-lite/0\": false,\n\"kawpow\": false\n},\n\"donate-level\": 0,\n\"donate-over-proxy\": 0,\n\"log-file\": null,\n\"pools\": [\n{\n\"algo\": null,\n\"coin\": \"monero\",\n\"url\": \"xmr-eu2.nanopool.org:14444\",\n\"user\":\n\"49dnvYkWkZNPrDj3KF8fR1BHLBfiVArU6Hu61N9gtrZWgbRptntwht5JUrXX1ZeofwPwC6fXNxPZfGjNEChXt\n\"pass\": \"x\",\n\"rig-id\": null,\n\"nicehash\": false,\n\"keepalive\": false,\n\"enabled\": true,\n\"tls\": false,\n\"tls-fingerprint\": null,\n\"daemon\": false,\n\"socks5\": null,\n\"m-select\": null\n}, {\n\"algo\": null,\n\n```\n\n-----\n\n```\n coin : monero,\n\"url\": \"xmr.f2pool.com:13531\",\n\"user\":\n\"49dnvYkWkZNPrDj3KF8fR1BHLBfiVArU6Hu61N9gtrZWgbRptntwht5JUrXX1ZeofwPwC6fXNxPZfGjNEChXt\n\"pass\": \"x\",\n\"rig-id\": null,\n\"nicehash\": false,\n\"keepalive\": false,\n\"enabled\": true,\n\"tls\": false,\n\"tls-fingerprint\": null,\n\"daemon\": false,\n\"socks5\": null,\n\"self-select\": null\n}\n],\n\"retries\": 5,\n\"retry-pause\": 5,\n\"syslog\": false,\n\"user-agent\": null,\n\"verbose\": 0,\n\"watch\": false,\n\"pause-on-battery\": false\n}\n\n## IOCs\n\n```\nFiles:\n```\nLdr.ps1\nd8a0ad2486c9662c8ac8c4370c7986bdc02dc6bdb8155dcc537f10dd00252d1d\nSys.exe\n618a1d0e1e1f6d0a34a5ab6c08acf058ffe46938f53b3357f388df1c7909fc45\nSysvr013.exe\n80bc76202b75201c740793ea9cd33b31cc262ef01738b053e335ee5d07a5ba96\n[kthreaddi].exe\na2f3ecd329d2713855257bf922b8a092cbb1193327ba197351804275286df7dd\n\n```\nNetwork:\n```\nExploit Source IP 191[.]252[.]219[.]71\nhttp://194.145.227[.]21/poc.xml\nhttp://194.145.227[.]21/pocwin.xml\nhttp://194.145.227[.]21/ldr.ps1\n\n## Detection\n\n```\n**Sigma**\n\n_please note that you need to enable webserver logging on the WebLogic server_\n\n**Snort**\n\n\n-----\n\n```\nSERVER WEBAPP Oracle WebLogic Server command injection attempt\nhttps://www.snort.org/rule_docs/1-56202\n\n```\n**Suricata**\n```\nET INFO Executable Download from dotted-quad Host\nAV POLICY HTTP request for .exe file with no User-Agent\nET CURRENT_EVENTS Terse alphanumeric executable downloader high likelihood of being\nhostile\nET INFO Packed Executable Download\nET INFO Executable Retrieved With Minimal HTTP Headers - Potential Second Stage\nDownload\nET INFO SUSPICIOUS Dotted Quad Host MZ Response\nET POLICY Cryptocurrency Miner Checkin\nETPRO POLICY XMR CoinMiner Usage\nETPRO TROJAN CoinMiner Known Malicious Stratum Authline\n\n```\n**YARA**\n\n\n-----\n\n```\n/\nYARA Rule Set\nAuthor: The DFIR Report\nDate: 2021-06-03\nIdentifier: 3580 WebLogic XMRig Miner\nReference: https://thedfirreport.com\n*/\n/* Rule Set ----------------------------------------------------------------- */\nimport \"pe\"\nrule sysrv013 {\nmeta:\ndescription = \"files - file sysrv013.exe\"\nauthor = \"The DFIR Report\"\nreference = \"https://thedfirreport.com\"\ndate = \"2021-06-03\"\nhash1 = \"80bc76202b75201c740793ea9cd33b31cc262ef01738b053e335ee5d07a5ba96\"\nstrings:\n$s1 = \"eDumped\" fullword ascii\n$s2 = \"OGhZVFNVRms6\" fullword ascii /* base64 encoded string '8hYTSUFk:' */\n$s3 = \"sdumpsebslemW^\" fullword ascii\n$s4 = \":ERNEL32.DLLODe?\" fullword ascii\n$s5 = \"pircsek\" fullword ascii\n$s6 = \"IDE5NC4xNDUu\" fullword ascii /* base64 encoded string ' 194.145.' */\n$s7 = \"333444\" ascii /* reversed goodware string '444333' */\n$s8 = \"aHx8d2dldCAtLXVzZXItYW\" fullword ascii /* base64 encoded string 'h||wget -user-a' */\n$s9 = \"h -c {%\" fullword ascii\n$s10 = \"?$?)?2?9?{\" fullword ascii /* hex encoded string ')' */\n$s11 = \";5A;6D2!D\" fullword ascii /* hex encoded string 'Zm-' */\n$s12 = \"3GRAt\\\\5\" fullword ascii\n$s13 = \"Gorget{\" fullword ascii\n$s14 = \"7!7&7,727\" fullword ascii /* hex encoded string 'ww'' */\n$s15 = \"* {5C;\" fullword ascii\n$s16 = \"fsftp_fC\" fullword ascii\n$s17 = \"POSTulL\" fullword ascii\n$s18 = \"gogetv/\" fullword ascii\n$s19 = \"\\\\x86+.pdb\" fullword ascii\n$s20 = \"_DIRCWD?\" fullword ascii\ncondition:\nuint16(0) == 0x5a4d and filesize < 10000KB and\n( pe.imphash() == \"406f4cbdf82bde91761650ca44a3831a\" or 8 of them )\n}\n\n```\n**MITRE**\n\nExploit Public-Facing Application – T1190\nCommand and Scripting Interpreter: PowerShell – T1059.001\nCommand and Scripting Interpreter: Windows Command Shell – T1059.003\nScheduled Task/Job: Scheduled Task – T1053.005\nBoot or Logon Autostart Execution: Registry Run Keys / Startup Folder – T1547.001\n\n\n-----\n\nImpair Defenses: Disable or Modify System Firewall – T1562.004\nMasquerading: Masquerade Task or Service – T1036.004\nProcess Discovery – T1057\nSystem Network Connections Discovery – T1049\nIngress Tool Transfer – T1105\nApplication Layer Protocol: Web Protocols – T1071.001\nResource Hijacking – T1496\n\nInternal case 3580\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2021/2021-06-03 - WebLogic RCE Leads to XMRig.pdf"
    ],
    "report_names": [
        "2021-06-03 - WebLogic RCE Leads to XMRig.pdf"
    ],
    "threat_actors": [],
    "ts_created_at": 1673536000,
    "ts_updated_at": 1743041139,
    "ts_creation_date": 1653761033,
    "ts_modification_date": 1653761033,
    "files": {
        "pdf": "https://archive.orkl.eu/d75cf01f6361a7f6e109e408656a8940b1a21c95.pdf",
        "text": "https://archive.orkl.eu/d75cf01f6361a7f6e109e408656a8940b1a21c95.txt",
        "img": "https://archive.orkl.eu/d75cf01f6361a7f6e109e408656a8940b1a21c95.jpg"
    }
}