# 0-raw-data/ORKL/library/jsons_entries/1c4a639a-321f-49ad-8f53-25a3247de8e6.json
{
    "id": "1c4a639a-321f-49ad-8f53-25a3247de8e6",
    "created_at": "2023-02-02T02:09:05.669767Z",
    "updated_at": "2025-03-27T02:05:55.751944Z",
    "deleted_at": null,
    "sha1_hash": "4afcec10a8db5d6cfb6632580a78438d8f3d5f74",
    "title": "2023-01-12 - QakBot Malware Used Unpatched Vulnerability to Bypass Windows OS Security Feature",
    "authors": "",
    "file_creation_date": "2023-02-01T07:48:44Z",
    "file_modification_date": "2023-02-01T07:48:44Z",
    "file_size": 4719146,
    "plain_text": "# QakBot Malware Used Unpatched Vulnerability to Bypass Windows OS Security Feature\n\n**[blog.eclecticiq.com/qakbot-malware-used-unpatched-vulnerability-to-bypass-windows-os-security-feature](https://blog.eclecticiq.com/qakbot-malware-used-unpatched-vulnerability-to-bypass-windows-os-security-feature)**\n\n## Executive Summary\n\n\n-----\n\nThis paper investigates a recent QakBot phishing campaign s ability to evade Mark-of-theWeb (MoTW) security features, allowing for escape from the designated security zone and\nsuccessful installation of malicious software on victim device.. Key observations:\n\nEclecticIQ analysts investigated QakBot phishing campaigns switching to a Zero-Day\nVulnerability to evade Windows Mark of the Web (MoTW). QakBot may be able to\nincrease its infection success rate as a result of the switch to a zero-day exploit.\n\nThe threat actor distributes QakBot using phishing emails with a malicious URL inside.\n\nWhen a victim user clicks on the malicious URL, it starts to download an encrypted ZIP\nfolder that contains an ISO image. If the ISO image is opened by victim, it will mount\nitself on a disk and open another File Explorer window that contains the final QakBot\nLoader as a JavaScript format which can be executed by a simple user click.\n\nThe final QakBot Loader (WW.js) contains a malformed digital signature to evade the\nMark of the Web (MoTW) Security feature on Windows OS. · EclecticIQ analysts\nobserved use of zero-day vulnerabilities is increasing among non-nation state cyber\ncriminals.\n\nLiving off the Land Binaries (LOLBINS) like Regsvr32.exe (2) and WScript.exe (3) are\nactively abused to execute QakBot Malware.\n\n## What is Mark of The Web (MoTW)?\n\nMark of the Web (MoTW) is used by Windows as a security feature across its product suite.\nThis feature works by checking downloaded executable files against a file whitelist that are\ndownloaded by Windows users.If the file is not on that list, Windows Defender SmartScreen\nwill show a warning message like image below and it will not execute the malware:\n\nFigure 1 – Windows SmartScreen warning\n\n\n-----\n\nThe MS Office Protected view feature is used to protect MS Office users against potential\nmalware in documents. Most of the MS Office file types flagged with MOTW will be opened\nwith PROTECTED VIEW:\n\nFigure 2 -MS Office document opened as Protected View\n\nMS Office is able to block macro enabled office document downloaded from the internet, if\nthe appropriate setting is enabled. Macros in MS Office files flagged with MOTW are disabled\nand a warning message is displayed to the user:\n\nFigure 3 – Macros blocked on downloaded Excel document\n\nWhen a Windows OS user downloads a file from the internet, it creates an Alternative Data\nStream (ADS) named Zone.Identifier and adds a ZoneId to this ADS in order to indicate the\nzone from which the file originates. This is a proactive security feature to prevent\ndownloading malicious files on untrusted source. Many Windows security features such as\nMicrosoft Office Protected view, SmartScreen, Smart App Control, and warning dialogs rely\non the presence of the MoTW to function correctly.\n\nAs the example image shows, details of MoTW alternate data streams on downloaded file\nfrom VirusTotal.ZoneID being used to identify a file, for example The following ZoneId values\nmay be used in a Zone.Identifier ADS:\n\n1. Local computer\n2. Local intranet\n3. Trusted sites\n4. Internet\n5. Restricted sites\n\n\n-----\n\nFigure 4 – Extracting ZoneID ADS on downloaded file\n\n## QakBot Campaign Observed Evading Windows Mark of the Web (MoTW)\n\nAt the beginning of November 2022, EclecticIQ analysts examined a recent campaign that\ndelivers QakBot (also called Qbot) to victim devices via phishing emails, executes by\nabusing multiple Living Off the Land Binaries (LOLBAS) and evades the Mark of the Web\n(MoTW) flag to increase the infection rate. Qakbot has been observed as an initial access\npoint for ransomware groups (4).\n\nThreat actors have used QakBot since 2007 (5) as a Banking Trojan to steal credit card\ninformation from victim devices. It evolved as initial access malware for remotely delivering\nadditional malicious payloads. Black Basta Ransomware gang used QakBot to create an\ninitial access point of victim's device and move laterally within an organization's network to\nexecute ransomware at the end of the kill chain.\n\nQakBot’s execution process is highlighted below:\n\n\n-----\n\nFigure 5 - QakBot Execution Flow\n\n## First Stage: Phishing Emails Containing Malicious URLs Deliver Qakbot Loader\n\nThe attack starts with a phishing email containing a malicious URL and ZIP password for\ndelivering the QakBot malware. Victims clicking on the URL download an encrypted ZIP\nfolder which can be unzipped with a password provided by attackers via phishing email. That\nunzipped file contains a randomly named malicious ISO image. The ISO image contains a\nfinal QakBot loader in form of a JavaScript file (WW.js) which is used to execute QakBot DLL\nin-memory of wermgr.exe (a Windows error reporting process).\n\n\n-----\n\nFigure 6 - Example of Phishing Email delivers QakBot Malware\n\n## Second Stage 2.1: In-Memory Execution of QakBot Malware via JavaScript Loader\n\nThe QakBot Loader can be executed by one of the most widely abused Living Off the Land\nBinaries And Scripts (LOLBAS) called wscript.exe (3). Threat Actors often abuse Windows\nbuilt in features to avoid detection. On Windows OS, JavaScript file extension can be\nexecuted by user click, upon the execution it uses Windows built in software called\nwscript.exe (3).\n\nFigure 7 - QakBot loader inside mounted ISO image.\n\nQakBot Loader deploys the Regsvr32.exe (2) command line tool as an obfuscated string to\nevade antivirus detections. When a user clicks on the WW.js, it will use Regsvr32.exe (2) to\nload the QakBot DLL, which is located under the port directory and is named\n\n\n-----\n\nresemblance.tmp.\n\nFigure 8 - QakBot Loader with malformed digital signature.\n\nFigure 9 - Resemblance.tmp contains MZ magic header which marking it executable.\n\n\n-----\n\nFigure 10 - Extracted malformed digital signature from JavaScript QakBot Loader\n\n## Second Stage 2.2: QakBot Loader uses Malformed Digital Signature to Evade Mark of the Web (MoTW)\n\nOn November 3rd, researcher Will Dormann (6) identified three different MoTW bypass\nmethods for bypassing the MoTW feature. On November, 8th, Microsoft released patches\n(CVE-2022-41049, CVE-2022-41091) (7) addressing two of the methods. The 3rd method using malformed digital signatures (CVE-2022-44698) (8) - patched on December 13 and is\nactively exploited in the wild.\n\nNormally, after executing the QakBot loader, Windows will display a warning message (see\nFigure 11) to avoid the execution. Because of the malformed digital signature, the loader\nbypasses the Mark of the Web (MoTW) flag, and the execution is proceeds without a\nWindows warning pop-up message.\n\n\n-----\n\nFigure 11 - Mark of the Web (MoTW) in action\n\n\n-----\n\nFigure 12 - Downloaded JavaScript file from untrusted URL automatically flagged by MoTW.\n\n## Third Stage: QakBot Uses Multiple Techniques to Evade Anti- Malware Scanners\n\nIn the next stage of the attack, QakBot injects itself inside the legitimate Windows Error\nReporting process (wermgr.exe) to evade behavior based anti-malware solutions.\n\n\n-----\n\nFigure 13 - Injected QakBot DLL\n\nMore information about the Living Off the Land Binaries Regsvr32.exe and WScript.exe can\nbe found via the links below.\n\nRegsvr32.exe (2)\nWScript.exe (3)\n\nFigure 14 - Process injection on wermgr.exe and LOLBAS observed in process tree.\n\nQakBot uses Windows API Hashing (Dynamic API Resolution) to evade signature-based\nanti-malware scanners. It hides the content of the import address table by XOR Encrypted\nAPI Hashing Algorithm called CRC32.\n\nBelow pictures showing Decompiled functions being used to perform API Hashing:\n\n\n-----\n\nFigure 15 - XOR Encrypted API Hashing.\n\nEclecticIQ analysts extracted the XOR key which is used to decrypt the content of APIs\nduring the execution time and used this key to decrypt other APIs for further analysis.\n\nFigure 16 - XOR Encryption key stored as static to decrypt the API hash.\n\nQakBot also uses the XOR encryption algorithm to hide its strings for minimizing AV\ndetection. Figure 10 shows encrypted strings are stored in the .rdata Section. They are\ndecrypted during run time.\n\n\n-----\n\nFigure 17 - XOR Encrypted strings hidden inside rdata section\n\nEclecticIQ analysts successfully decrypted the XOR encrypted strings used by QakBot. The\ndecrypted strings are used by QakBot for testing the internet connection of the victim device,\nconducting a sandbox check, gaining persistence on the victim device by abusing Schedule\nTask, and gathering victim computer information upon the attacker’s request through a\ncommand-and-control (C2) server.\n\nFigure 18 – Decrypted Strings from QakBot Malware\n\n\n-----\n\n## Fourth Stage: Command and Control (C2) Connection\n\nAfter successful execution, QakBot checks its internet connectivity and will send multiple\nPOST requests to its C2 servers.\n\nQakBot checks internet availability on victim's device:\n\nFigure 19 - QakBot malware checking Internet availability\n\nC2 protocol uses JSON object encapsulation with a RC4 Encrypted message which is\nencoded with Base64.\n\nFigure 20 - QakBot performs command and control connections\n\nRaw example of an HTTP POST request sent by QakBot to its C2:\n\nMITRE ATT&CK\n\n**Technique Name** **TTP ID**\n\nUser Execution: Malicious Link T1204.001\n\nSystem Binary Proxy Execution: Regsvr32 T1218.010\n\nCommand and Scripting Interpreter: JavaScript T1059.007\n\n\n-----\n\nPhishing: Spearphishing Link T1566.002\n\nApplication Layer Protocol: Web Protocols T1071.001\n\nProcess Injection: Process Hollowing T1055.012\n\nObfuscated Files or Information T1027\n\nObfuscated Files or Information: Dynamic API Resolution T1027.007\n\nSystem Information Discovery T1082\n\nScheduled Task/Job: Scheduled Task T1053.005\n\nVirtualization/Sandbox Evasion: System Checks T1497.001\n\nWindows Management Instrumentation T1047\n\n## Indicators:\n\n**FIle Name** **SHA 256 Hash**\n```\n resemblance.tm\n p 8ca16991684f7384c12b6622b8d1bcd23bc27f186f499c2059770ddd3031\n           f274\n UY76.img\n           26f5bc698dfec8e771b781dc19941e2d657eb87fe8669e1f75d9e5a1bb4d\n           b1db\n WW.js\n           c5df8f8328103380943d8ead5345ca9fe8a9d495634db53cf9ea3266e353\n           a3b1\n Injected QakBot-dll 6fb41b33304b65e6e35f04e8cc70f7a24cd36e29bbb97266de68afcf113f\n           9a5f\n\n```\nFind the data for [COMMAND AND CONTROL SERVER C2](https://www.eclecticiq.com/hubfs/_blogs/corporate-blog/2023/command-control-qakbot.txt)\n\nFind the data for [YARA RULES](https://www.eclecticiq.com/hubfs/_blogs/corporate-blog/2023/yara-rules-qakbot-final.txt)\n\n## About EclecticIQ Intelligence & Research Team\n\n\n-----\n\nEclecticIQ is a global provider of threat intelligence, hunting, and response technology and\n[services. Headquartered in Amsterdam, the EclecticIQ Intelligence & Research Team is](https://www.eclecticiq.com/eclecticiq-intelligence-research-team)\nmade up of experts from Europe and the U.S. with decades of experience in cyber security\nand intelligence in industry and government.\n\nWe would love to hear from you. Please send us your feedback by emailing us\nat [research@eclecticiq.com or fill in the EclecticIQ Audience Interest Survey to drive our](mailto:research@eclecticiq.com)\nresearch towards your priority area.\n\n## Structured Data\n\nFind the Analyst Prompt and earlier editions in our public TAXII collection for easy use in\nyour security stack.\n\n[TAXII v1 Discovery services: https://cti.eclecticiq.com/taxii/discovery](https://cti.eclecticiq.com/taxii/discovery)\n\nPlease refer to our [support page for guidance on how to access the feeds.](https://www.eclecticiq.com/public-feed-manual)\n\n## You might also be interested in:\n\n[Network Environment-Focused Conversations Needed in Approaches to Cyber Security](https://blog.eclecticiq.com/emotet-downloader-document-uses-regsvr32-for-execution)\n\n[Emotet Downloader Document Uses Regsvr32 for Execution](https://blog.eclecticiq.com/the-analyst-prompt-10-ai-facial-recognition-used-in-ukraine/russia-war-prone-to-vulnerabilities)\n\n[AI Facial Recognition Used in Ukraine/Russia War Prone to Vulnerabilities](https://blog.eclecticiq.com/the-analyst-prompt-10-ai-facial-recognition-used-in-ukraine/russia-war-prone-to-vulnerabilities)\n\n## Appendix\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "4e039bf8-7239-4999-8ec0-eb202d1465af",
            "created_at": "2022-12-04T11:19:36.52593Z",
            "updated_at": "2022-12-04T11:19:36.52593Z",
            "deleted_at": null,
            "name": "ORKL",
            "url": "https://github.com/ORKL/library",
            "description": "ORKL Community Contributed Content",
            "reports": null
        },
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2023/2023-01-12 - QakBot Malware Used Unpatched Vulnerability to Bypass Windows OS Security Feature.pdf"
    ],
    "report_names": [
        "2023-01-12 - QakBot Malware Used Unpatched Vulnerability to Bypass Windows OS Security Feature.pdf"
    ],
    "threat_actors": [],
    "ts_created_at": 1675303745,
    "ts_updated_at": 1743041155,
    "ts_creation_date": 1675237724,
    "ts_modification_date": 1675237724,
    "files": {
        "pdf": "https://archive.orkl.eu/4afcec10a8db5d6cfb6632580a78438d8f3d5f74.pdf",
        "text": "https://archive.orkl.eu/4afcec10a8db5d6cfb6632580a78438d8f3d5f74.txt",
        "img": "https://archive.orkl.eu/4afcec10a8db5d6cfb6632580a78438d8f3d5f74.jpg"
    }
}