# 0-raw-data/ORKL/library/jsons_entries/368f449a-c527-41c8-a00a-68b135e71c4a.json
{
    "id": "368f449a-c527-41c8-a00a-68b135e71c4a",
    "created_at": "2023-01-12T15:10:34.876342Z",
    "updated_at": "2025-03-27T02:09:18.429667Z",
    "deleted_at": null,
    "sha1_hash": "0d480cc82a617a0ec7c736e8f2e5e39bd73978e5",
    "title": "2021-09-29 - Zloader Campaigns at a Glance",
    "authors": "",
    "file_creation_date": "2022-05-28T16:01:40Z",
    "file_modification_date": "2022-05-28T16:01:40Z",
    "file_size": 2270609,
    "plain_text": "# Zloader Campaigns at a Glance\n\n**[trendmicro.com/vinfo/us/security/news/cybercrime-and-digital-threats/zloader-campaigns-at-a-glance](https://www.trendmicro.com/vinfo/us/security/news/cybercrime-and-digital-threats/zloader-campaigns-at-a-glance)**\n\n\n-----\n\nView infographic: Zloader Campaigns at a Glance\n\nThe ZBOT (aka Zeus) trojan has been one of the most prolific and enduring malware families\nof the past 20 years. After its first [appearance in 2006, its source code was leaked in 2011,](https://www.trendmicro.com/vinfo/us/security/news/cybercrime-and-digital-threats/online-banking-trojan-brief-history-of-notable-online-banking-trojans)\nleading to a plethora of new variants that plagued organizations over the succeeding years.\n\nOne of the most notable recent ZBOT variants is Zloader. First compiled under the name\nSilent Night in late 2019, it has evolved from being an information stealer to a multipurpose\ndropper that provides malicious actors the means to install and execute other malware and\ntools such as Cobalt Strike, DarkSide, and Ryuk. In addition, it has other capabilities, such\nas the ability to provide remote access to attackers and install plug-ins for additional routines.\n\nZloader has multiple delivery methods, such as via email campaigns or downloads by other\nmalware and hacking tools. One of the most basic yet reliable methods for individuals and\norganizations to avoid being infected by Zloader and other malware with similar arrival\ntechniques is to apply security best practices to their emails. This includes avoiding\ndownloading attachments or selecting links from emails that look suspicious or appear to be\nout of context.\n\nZloader’s versatility has made it a popular and effective campaign tool for any threat actor\nthat is willing to pay for it. We already witnessed this in past campaigns — some of which\ntook advantage of current events such as the Covid-19 pandemic — and we can expect to\nsee it again in future campaigns from other threat actors.\n\nOrganizations can mitigate the impact of Zloader by employing robust security solutions and\nservices. Trend Micro’s robust native XDR capabilities are tied together by Trend Micro\nVision One™which connects email endpoints servers cloud workloads and networks in\n\n\n-----\n\norder to provide a better context and perspective of the entire chain of events of an attack,\nwhile also allowing security personnel to investigate and act from a single place.\n\n[Furthermore, managed security services, such as Trend Micro™ Managed XDR, provides](https://www.trendmicro.com/en_us/business/services/managed-xdr.html)\nexpert threat monitoring, correlation, and analysis from experienced cybersecurity\nprofessionals via a single and capable source of detection, analysis, and response. This\nexpertise is further bolstered by AI-optimized, Trend Micro solutions that draw from global\nthreat intelligence.\n\n## MITRE ATT&CK techniques\n\nZloader uses the following tactics and techniques, as mapped out according to the MITRE\nATT&CK Matrix.\n\n\nTactic MITRE ID and\nTechnique\n\nInitial Access T1189 - Drive-by\nCompromise\n\n\nDetails\n\nZloader can be downloaded through drive-by\ncompromise via Malsmoke, RIG Exploit Kit, and\nSpelevo\n\n\nT1566 Phishing\n\n\nZloader can arrive via\nphishing emails with\nattached XLS\ndownloader files\n\n\nExecution T1204 - User Execution User can execute the XLS Zloader downloader\nfile manually\n\n\n-----\n\nT1064 Scripting\n\nT1059 Command\nand Scripting\nInterpreter\n\nT1106 Native API\n\n\nZloader can be\ndownloaded by VBS or\nJavascripts\n\nZloader hooks native\nAPI from user32.dll and\nntdll.dll to redirect\nexecution to Zloader\nDLL\n\n\nPersistence T1060 - Registry Run\nKeys/Startup Folder\n\nT1547- Boot\nor Logon\nAutostart\nExecution\n\n\nCreates persistence using the following\nregistry: HKCU\\SOFTWARE\\Microsoft\\\nWindows\\CurrentVersion\\Run\n\nZloader injects its loader or core component to\nmsiexec.exe\n\n\nPrivilege\nEscalation\n\n\nT1055 - Process\nInjection\n\n\n-----\n\nDefense\nEvasion\n\nT1140 –\nDeobfuscate/\nDecode Files\nor\nInformation\n\nT1497 Virtualization/\nSandbox\nEvasion\n\nCredential\nAccess\n\nT1539 Steal Web\nSession\nCookie\n\n\nT1027 - Obfuscated files\nor information\n\nZloader performs XOR\nto decode obfuscated\nstrings and information\n\nZloader downloader\nscripts check if it is\nrunning in a virtual\nenvironment and will not\nexecute properly if it is\n\n\nInstead of presenting arithmetic functions in a\nstandardized manner and directly hardcoding\nconstants, Zloader tries to confuse the analyst\nby obfuscating these in a form of various,\ndedicated functions\n\n\nT1056 - Input Capture Zloader captures keystrokes on browsers\n\nZloader steals cookies\nfrom Chrome, Firefox,\nand Internet Explorer\n\n\n-----\n\nDiscovery T1083 - File and\nDirectory Discovery\n\nT1012 Query\nRegistry\n\nCollection T1185 - Man in the\nBrowser\n\nT1179 Hooking\n\n\nZloader steals cookies by discovering files from\nspecific directories like \\Mozilla\\Firefox\\Profiles\n\nZloader has to install its own (fake) certificate,\nand has to run a local proxy before deploying a\nMan-In-TheBrowser (MITB) attack\n\n\n-----\n\nC2 is encrypted via RC4 and XORing algorithm\nwhere each character of the string is XORed\nwith the preceding character which was already\nXORed\n\n\nCommand &\nControl\n\n\nT1001 - Data\nObfuscation\n\n\n-----\n\nT1090 Proxy\n\nT1071Application\nLayer\nProtocol\n\nT1219 Remote\nAccess\nSoftware\n\n\nZloader components\ninjected into browsers\nare responsible for\nredirecting traffic via\nproxy\n\nThe following\ncommands are\naccepted:\n\nuser_execute download an executable\ninto the %TEMP% folder\nand run it (optionally\nwith parameters)\n\nuser_cookies_get - steal\ncookies from all known\nbrowsers\n\nuser_url_block - block\nURL access for the\ncurrent user\n\nbot_uninstall - complete\nremoval of the bot from\nthe current user\n\nuser_password_get –\nsteal passwords from\ntargeted browsers\n\nuser_files_get – search\nand upload documents\nof the victims (.txt,\ndocx,, .xls, .wallet.dat)\n\nZloader downloads and\nexecutes VNC tool to\ncontrol victim machine\n\n\nExfiltration T1041 - Exfiltration Over\nC&C Channel\n\n## Indicators of Compromise\n\n\nData collected by Zloader, such as stolen\ncookies, screenshot, and process list, are\nexfiltrated to C&C server\n\n\n[The IOCs for Zloader can be found in this appendix.](https://documents.trendmicro.com/assets/txt/IOCs-zloader-campaigns-at-a-glance.txt)\n\n\n-----\n\nHIDE\n\n\n-----\n\n-----\n\n-----\n\n**Like it? Add this infographic to your site:**\n\n1. Click on the box below.  2. Press Ctrl+A to select all.  3. Press Ctrl+C to copy.  4. Paste\nthe code into your page (Ctrl+V).\n\nImage will appear the same size as you see above.\n\n[Posted in Cybercrime & Digital Threats,](https://www.trendmicro.com/vinfo/us/security/news/cybercrime-and-digital-threats) [Infographics,](https://www.trendmicro.com/vinfo/us/security/news/infographics) [Targeted Attacks](https://www.trendmicro.com/vinfo/us/security/news/targeted-attacks)\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2021/2021-09-29 - Zloader Campaigns at a Glance.pdf"
    ],
    "report_names": [
        "2021-09-29 - Zloader Campaigns at a Glance.pdf"
    ],
    "threat_actors": [
        {
            "id": "610a7295-3139-4f34-8cec-b3da40add480",
            "created_at": "2023-01-06T13:46:38.608142Z",
            "updated_at": "2025-03-27T02:00:02.87217Z",
            "deleted_at": null,
            "main_name": "Cobalt",
            "aliases": [
                "Cobalt Gang",
                "GOLD KINGSWOOD",
                "COBALT SPIDER",
                "G0080",
                "Mule Libra",
                "Cobalt Group"
            ],
            "source_name": "MISPGALAXY:Cobalt",
            "tools": [],
            "source_id": "MISPGALAXY",
            "reports": null
        }
    ],
    "ts_created_at": 1673536234,
    "ts_updated_at": 1743041358,
    "ts_creation_date": 1653753700,
    "ts_modification_date": 1653753700,
    "files": {
        "pdf": "https://archive.orkl.eu/0d480cc82a617a0ec7c736e8f2e5e39bd73978e5.pdf",
        "text": "https://archive.orkl.eu/0d480cc82a617a0ec7c736e8f2e5e39bd73978e5.txt",
        "img": "https://archive.orkl.eu/0d480cc82a617a0ec7c736e8f2e5e39bd73978e5.jpg"
    }
}