# 0-raw-data/ORKL/library/jsons_entries/38cb45cd-3d04-480a-9414-f1f484187da7.json
{
    "id": "38cb45cd-3d04-480a-9414-f1f484187da7",
    "created_at": "2023-01-12T15:06:14.977414Z",
    "updated_at": "2025-03-27T02:09:18.367494Z",
    "deleted_at": null,
    "sha1_hash": "debc0a337dc20447e68e7abd8966bc15cb4b13ee",
    "title": "2020-11-12 - Egregor – Prolock- Fraternal Twins -",
    "authors": "",
    "file_creation_date": "2022-05-27T21:32:23Z",
    "file_modification_date": "2022-05-27T21:32:23Z",
    "file_size": 883813,
    "plain_text": "# Egregor – Prolock: Fraternal Twins ?\n\n**intrinsec.com/egregor-prolock/**\n\nEquipe CERT November 12, 2020\n\n_CERT Intrinsec has faced since the beginning of September several cases involving Egregor and Prolock ransomwares. This article aims at_\n_presenting Egregor and Prolock techniques, tactics and procedures, as well as sharing indicators of compromise and highlighting actions of_\n_the threat actor operating both ransomwares, according to collected intelligence and TTPs analysis._\n\n_On one hand, Egregor has a similar strategy to other ransomwares, as it exfiltrates data, encrypts files and publishes them on its website in_\n_order to make victims pay the ransom. It is active since the beginning of September 2020 and impacts many sectors from insurance to_\n_transport. Its goal is lucrative._\n\n_On the other hand, Prolock, successor of PwndLocker, is active since March 2020. As many other ransomwares, it targets big companies with_\n_ransoms going from 35 to 255 bitcoins (400 000 to 3 000 000 $). Its goal is, as far as we know, only lucrative. Prolock is active mainly in_\n_Northern America and in Europe and impacts several sectors, such as health, construction, finance, legal, etc._\n\n## Kill chain\n\n Egregor ransomware analysis\n\n\n-----\n\nWe were not able to find identified specific initial accesses, we found traces of Qakbot during investigations but we could not identify how it\nwas dropped on information systems. We observed many potential intrusion vectors on patient 0 (many malwares were found on the machine).\n\n**Internal reconnaissance**\n\nPrior to privilege escalation, Egregor proceeds to Active Directory reconnaissance using tools such as Sharphound or AdFind. These tools are\nused to gather information about users, groups, computers, and so on. They aim as well at finding the best compromission paths.\n\n**Privilege escalation**\n\nDuring investigations, Egregor compromises Active Directory in order to become domain admin.\n\n**Lateral movement**\n\nEgregor moves laterally on information systems using CobaltStrike SMB beacons. This feature allows an attacker to use SMB named pipes\n(logical connections between a client and a server) to communicate commands through the information system revealing C2 IP address.\n\nThe following command line is a service created by CobaltStrike and can be found in Windows Event Logs (event id 7045). It runs an encoded\npowershell command.\n\n_CobaltStrike service execution_\n[It is possible to deobfuscate CobaltStrike payloads (base64, gunzip and XOR operations) using CyberChef[1]:](https://github.com/mattnotmax/cyberchef-recipes#recipe-28---de-obfuscation-of-cobalt-strike-beacon-using-conditional-jumps-to-obtain-shellcode)\n\n_CobaltStrike payload deobfuscation_\n\n**C2 Communication**\n\nOnce settles on victim’s information systems, Egregor communicates with its Command and Control servers via HTTPS protocol so as to drop\nscripts or dynamic link libraries on infected hosts. You can find the list of C2 identified during investigations in section “IP Addresses”.\n\n**Data exfiltration**\n\nEgregor masquerades svchost.exe process to launch an RClone client in order to exfiltrate data. RClone aims at managing files in cloud, it\ndeals with multiple systems and protocols. The RClone configuration file, in plain text, is dropped by the attacker with the binary. Based on\ninvestigations and OSINT, we know that Egregor used at least three different configurations to exfiltrate data.\n\n\n-----\n\n_RClone Configuration File (WebDav)_\n\n_RClone Configuration File (sFTP)_\n\n\n_RClone_\n\n_Configuration File (DropBox)_\n\n**Defense evasion**\n\nTo evade protections, Egregor create a Group Policy Object to disable Windows Defender and try to takedown any anti-virus console prior to\nransomware execution:\n```\nDisplay name: New Group Policy Object\nVersion: 1\nregistry.pol content:\n- Key path: Software\\Policies\\Microsoft\\Windows Defender\n- Data name: DisableAntiSpyware\n- Value type: 0x04 (REG_DWORD)\n- Data value: 0x01\n\n```\n**Ransomware execution**\n\nEgregor downloads custom dynamic link libraries (b.dll, q.dll, etc) using bitsadmin and execute them on victim’s systems to encrypt data.\n\n_DLL download and execution_\n\n## Prolock ransomware analysis\n\n**Initial Access**\n\nOne of the intrusion vectors is malspam. Indeed, Emotet is used to initiate infection on several user workstations and to drop Qakbot. Emotet\nused legitimate documents after taking control of some user’s email accounts. These documents contain a payload which tries to download a\nbinary file from different URLs, as following. On infected systems, after the execution of the binary retrieved by Emotet code, few files are\ncreated (typical Qakbot operation):\n\n\\AppData\\Roaming\\Microsoft\\Jfayae\\vatgrcxt.exe\n\\koyxogldypnalvvlyxpw.exe\n\\AdFind.exe\n\n_Code embedded in the malicious document_\n\n\n-----\n\n_Powershell payload (decoded from base64)_\n\n_Deobfuscated powershell code (used to download F889k6.exe)_\nUnfortunately, we were not able to retrieve F889k6.exe neither from compromised systems nor from URLs, which were already down by that\ntime.\n\n**Internal reconnaissance**\n\nProlock proceeds to Active Directory reconnaissance AdFind tool to gather information about users, groups, computers so as to prepare\nexfiltration and ransomware execution.\n\n**Privilege escalation**\n\nDuring investigations, Prolock compromises Active Directory in order to become domain admin.\n\n**Lateral movement**\n\nProlock uses batch scripts to enable RDP on targeted hosts. We found the script below during one of our cases. The same script has already\nbe found on Prolock cases.\n\nFollowing actions are performed by the script:\n\nEnable Remote Desktop connections by setting fDenyConnections to 0.\nStart Microsoft Protection Service.\nSet a rule in Windows firewall to activate RDP service.\nModify RDP-Tcp registry key.\n\n_rdp.bat script (enable RDP connections)_\n\n**Data exfiltration**\n\nWe did not see any use of RClone during incident responses involving Prolock.\n\n\n-----\n\nProlock uses different scripts and files to encrypt victim’s data. It retrieves all these files from 185.238.0[.]233, the latter hosting as well\n**Egregor dynamic link libraries. The first script wmi_md.bat (wmi_u.bat works the same way) proceeds the following actions on each host**\nwhose IP address is in the file list_md.txt (or list_u.txt):\n\nConnect to the host using a compromised account\nDrop connect.bat and office.txt on the host\nExecute connect.bat using WMI command-line\nWrite host IP address in log.dat file\nCancel the network connection\n\n_Script deploying ransomware on information system (wmi_md.bat)_\nIn addition, we found a script that uses bitsadmin to download office.txt and connect.bat from 185.238.0[.]233.\n\n_Code from eb1.bat_\nThe script connect.bat contains the following encoded powershell payload.\n\n_Powershell payload from connect.bat_\nAfter decoding and deobfuscating it, we got to know that it is used to load office.txt in memory and execute it.\n\n_Decoded and deobfuscated payload_\n_Office.txt analysis is not yet complete, but we believe that it is the ransomware, based on system events correlation._\n\n## Relations between Egregor & Prolock\n\nDuring recent investigations, we observed common indicators of compromise and techniques between Egregor and Prolock. These common\npoints are presented below:\n\nThe IP address 185.238.0[.]233 hosts both Egregor’s dynamic link libraries and Prolock files (especially scripts used to run the\nransomware). You can find more information about Prolock TTP in the next section.\nBoth WIN-799RI0TSTOF and WIN-4K804V6ADVQ hostnames of potential VPS have been seen during Prolock and Egregor cases.\n_list_md.txt and list_u.txt files were involved in both Egregor and Prolock cases (of course, their content depends on the victim’s_\ninformation system).\nThe use of bitsadmin in eb*.bat scripts to download dll (Egregor) or scripts (Prolock) is another common point between these threat\nactors.\n\n\n-----\n\n_d e e b a y_ as bee see o bot g ego a d o oc cases\n[Even if we did not notice exfiltration using RClone in our Prolock cases, we know that this threat actor uses it[2].](https://www.group-ib.com/blog/prolock)\n\n_Timeline of incident responses involving Prolock and Egregor insisting on common indicators of compromise_\n\n## Indicators of compromise\n\n**Incident response**\n\nBinaries\n\n\n**File** **Size**\n**(bytes)**\n\n**md.exe** 4516928\n4183104\n\n\n**MD5** **SHA1** **SHA256**\n\n\n5ed9fb5fc74c6fdb3537629e9b23437a\n\nN/A\n\n\n67424175620be87fd3b2810ba5eba0d9e0bee49f\n\n7e0018e18f6bd230366a2b6f031c52ee8899f8dc\n\n\nfec51f04710e3aa\nN/A\n\n\n**svchost.exe** 42043904 4a97c4345aabf9dd922d29687c95ac66 f54bf6a4c6f7c3d0077d152a094e3c7738cf0bd1 5bc506b9f61ecec\n\n**main_target1.exe** 4516416 a3e1ea9438e293ec8fae62c39ea3f0e4 e9581cb5161f10f5e99e0cb6c30c201e6e844676 089bb9d18b3faf4\n\n**b.dll** 808960 a654b3a37c27810db180822b72ad6d3e d2d9484276a208641517a2273d96f34de1394b8e 4c9e3ffda0e6632\n\n**q.dll** 784896 520ee511034717f5499fb66f9c0b76a5 3a33de9a84bbc76161895178e3d13bcd28f7d8fe a5989c480ec650\n\n**qymrkrr.exe** N/A N/A N/A N/A\n\n**fxmgwk.exe** N/A N/A N/A N/A\n\n**cthwiilhz.exe** N/A N/A N/A N/A\n\n**erkftj.exe** N/A N/A N/A N/A\n\n**c6d7790.exe** N/A N/A N/A N/A\n\n**c6d7790.exe** N/A N/A N/A N/A\n\n**a31b29b.exe** N/A N/A N/A N/A\n\n**ed53e67.exe** N/A N/A N/A N/A\n\n**3f2eb85.exe** N/A N/A N/A N/A\n\nScripts\n\n\n**File** **Size**\n**(bytes)**\n\n**e.bat** 156\n\n157\n\n\n**MD5** **SHA1** **SHA256**\n\n\n7375083934dd17f0532da3bd6770ab25\nN/A\n\n\nac6d919b313bbb18624d26745121fca3e4ae0fd3\n1be22505a25f14fff1e116fafcaae9452be325b1\n\n\nf0adfd3f89c9268953f93bfd\nN/A\n\n\n**eb.bat** 58 N/A 9dacb159779d5e57798632bac74ae5b880cf1ec8 N/A\n\n**eb1.bat** 253 3872e7caaede9ee1ce8f37435dcaf836 8f166dfeb2fd8780de0e3dbdb25d0fdb373f58de c9df055f380100a7301f525\n\n**connect.bat** 7004 6cebf3c01844520e8b27023d8f47a0ed f5b14cc494303c91456bb50e7816358b6766a5b8 bd1dba49596c04677a2aa\n\n\n-----\n\n**wmi_u.bat** 421 9deca294973f6d52f9506240b104079c f098e6931eb32f9d28f681ad6fd2716a65b7f140 87a699923f3edeb6ce631f\n\n**wmi_md.bat** 416 463d45502447c7aa58538159eccc1a1a 4bad78fccfc69f4f9ac619dd9a8a9f70c3cc3ed0 a9d3c1d779550b0031285\n\n**rdp.bat** 329 dc1aafc01b5068eef6c2ed4cfd6864ed eb43350337138f2a77593c79cee1439217d02957 ac49c114ef137cc198786a\n\nOther files\n\n\n**File** **Size**\n**(bytes)**\n\n\n**SHA1**\n\n\n**svchost.conf** 155 bae4323aa7fa3e4de9ab021d72ecd84de795351b\n\n**office.txt** 30608 4769a775fd4a2c29b433736a59dc4277354a54f2\n\n**list_u.txt** 4560 9269 3b59fdff922497dc24d7cec0b219e93334e81221\n33c776f25ed3bb6011bfe96c13467815fb993289\n\n**list_md.txt** 4773 9633 1a3c149a2720f001a0a475ae978114090f3ed720\naaf4374400c63b0dae41f67bd90cd2ebb2c159db\n\n**list3.txt** 4560 3b59fdff922497dc24d7cec0b219e93334e81221\n\n\n\n**[HOW TO RECOVER**\n**FILES].TXT**\n\n\n1085 620311402640b1547d59722b63f19fab082a57af\n\n\n**RECOVER-FILES.txt** N/A N/A\n\nC2 Domain names\n\namajai-technologies[.]network\namajai-technologies[.]industries\n\nIP adresses\n\n_IP addresses hosting dll and scripts:_\n\n185.238.0[.]233\n45.153.242[.]129\n_Server using for potential data exfiltration_\n\n93.190.140[.]75\n_IP addresses communicating with infected systems through CobaltStrike_\n\n23.254.229[.]82\n192.236.209[.]151\n_Potential VPS IP addresses_\n\n217.138.219[.]138\n185.212.170[.]158\n23.106.215[.]67\n\nPotential VPS Hostnames\n\nWIN-799RI0TSTOF\nWIN-4K804V6ADVQ\nDESKTOP-LHC2KTF\nDESKTOP-93VHU8M\n\n**Threat Intelligence**\n\nUsing IOC collected during incident responses, we hunted some other Egregor files, especially from 185.238.0[.]233. We found similar\ndynamic link libraries (a.dll, p.dll, etc), as well as the RClone configuration file we presented in section “Data Exfiltration“.\n\n\n-----\n\n**File** **Size (bytes)** **SHA1**\n\n**b.dll** 808960 d2d9484276a208641517a2273d96f34de1394b8e\n\n**hnt.dll** 498688 38c88de0ece0451b0665f3616c02c2bad77a92a2\n\n**kk.dll** 498176 09d8c91ccefd699fb5ac1aaebeeebee25170fe1a\n\n**p.dll** 784896 8768cf56e12a81d838e270dca9b82d30c35d026e\n\n**p.dll** 500224 fafd32e972ebb33b187bfb1ebf1a6ecb1d2d7239\n\n**sed.dll** 806400 b7170443ea2b73bca3d16958712ee57cb4869d5b\n\nCobaltStrike C2 Domain names[3]\n\natakai-technologies[.]space\natakai-technologies[.]website\natakai-technologies[.]host\natakai-technologies[.]online\natakai-technologies[.]work\nakamai-technologies[.]host\nakamai-technologies[.]site\nakamai-technologies[.]space\nakamai-technologies[.]digital\nakamai-technologies[.]website\nakamai-technologies[.]online\namajai-technologies[.]host\namajai-technologies[.]website\namajai-technologies[.]network\namajai-technologies[.]digital\namajai-technologies[.]space\namajai-technologies[.]tech\namajai-technologies[.]industries\namamai-tecnologies[.]space\namamai-tecnologies[.]cloud\namamai-tecnologies[.]digital\n\n\n-----\n\na ata tec o og es[ ] ebs te\namatai-technologies[.]digital\namatai-technologies[.]space\namatai-technologies[.]site\n\n## MITRE ATT&CK\n\n**Prolock**\n\n**Tactic** **Technique**\n\n**Initial Access** Phishing (T1566): Spearphishing attachment (T1566.001)\n\n**Execution** User Execution (T1204): Malicious File (T1204.002) Windows Management Instrumentation (T1047)\n\n**Persistence** Scheduled Task/Job (T1053): Scheduled Task (T1053.005) Valid Accounts (T1078)\n\n**Discovery** Account Discovery (T1087) Domain Trust Discovery (T1482) Permission Groups Discovery (T1069): Domain\nGroups (T1069.001)\n\n**Lateral Movement** Remote Services (T1021): Remote Desktop Protocol (T1021.001) Valid Accounts (T1078)\n\n\n**Command and**\n**Control**\n\n\nIngress Tool Transfert (T1105)\n\n\n**Impact** Data encrypted for impact (T1486)\n\n**Egregor**\n\n**Tactic** **Technique**\n\n**Execution** Scheduled Task/Job (T1053): Scheduled Task (T1053.005) Services Execution (T1569): Service Execution (T1569.002)\nWindows Management Instrumentation (T1047)\n\n**Persistence** Create or modify system process (T1543): Windows Service (T1543.003)\n\n\n**Defense**\n**Evasion**\n\n\nImpair Defenses (T1562): Disable or modify tools (T1562.001)\n\n\n**Discovery** Account Discovery (T1087) Domain Trust Discovery (T1482) Permission Groups Discovery (T1069): Domain Groups\n(T1069.001)\n\n\n**Lateral**\n**Movement**\n\n**Command and**\n**Control**\n\n\nRemote Services (T1021): SMB/Windows Admin Shares (T1021.002)\n\nApplication Layer Protocol (T1071)\n\n\n**Exfiltration** Exfiltration over web service (T1567): Exfiltration to Cloud Storage (T1567.002)\n\n**Impact** Data encrypted for impact (T1486)\n\n## References\n\nCyberChef recipe to deobfuscate CobaltStrike payloads : [1] https://github.com/mattnotmax/cyberchef-recipes#recipe-28—de-obfuscation-ofcobalt-strike-beacon-using-conditional-jumps-to-obtain-shellcode\n\n[Article from Group-IB about Prolock : [2]https://www.group-ib.com/blog/prolock](https://www.group-ib.com/blog/prolock)\n\n[CobaltStrike C2 List :](https://www.intrinsec.com/wp-admin/post.php?post=221234&action=edit#_ftnref1) [[3]https://twitter.com/smoothimpact/status/1308033998371905538](https://twitter.com/smoothimpact/status/1308033998371905538)\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2020/2020-11-12 - Egregor – Prolock- Fraternal Twins -.pdf"
    ],
    "report_names": [
        "2020-11-12 - Egregor – Prolock- Fraternal Twins -.pdf"
    ],
    "threat_actors": [
        {
            "id": "610a7295-3139-4f34-8cec-b3da40add480",
            "created_at": "2023-01-06T13:46:38.608142Z",
            "updated_at": "2025-03-27T02:00:02.87217Z",
            "deleted_at": null,
            "main_name": "Cobalt",
            "aliases": [
                "Cobalt Gang",
                "GOLD KINGSWOOD",
                "COBALT SPIDER",
                "G0080",
                "Mule Libra",
                "Cobalt Group"
            ],
            "source_name": "MISPGALAXY:Cobalt",
            "tools": [],
            "source_id": "MISPGALAXY",
            "reports": null
        }
    ],
    "ts_created_at": 1673535974,
    "ts_updated_at": 1743041358,
    "ts_creation_date": 1653687143,
    "ts_modification_date": 1653687143,
    "files": {
        "pdf": "https://archive.orkl.eu/debc0a337dc20447e68e7abd8966bc15cb4b13ee.pdf",
        "text": "https://archive.orkl.eu/debc0a337dc20447e68e7abd8966bc15cb4b13ee.txt",
        "img": "https://archive.orkl.eu/debc0a337dc20447e68e7abd8966bc15cb4b13ee.jpg"
    }
}