# 0-raw-data/ORKL/library/jsons_entries/39a9e68d-916c-4154-a14a-36ed4ccb9492.json
{
    "id": "39a9e68d-916c-4154-a14a-36ed4ccb9492",
    "created_at": "2023-01-12T15:10:21.914632Z",
    "updated_at": "2025-03-27T02:15:41.794179Z",
    "deleted_at": null,
    "sha1_hash": "e0974628ddea51d61fa5c71f4b880039f3b71cb7",
    "title": "2020-06-03 - New LNK attack tied to Higaisa APT discovered",
    "authors": "",
    "file_creation_date": "2022-05-28T03:32:26Z",
    "file_modification_date": "2022-05-28T03:32:26Z",
    "file_size": 1435280,
    "plain_text": "# New LNK attack tied to Higaisa APT discovered\n\n**blog.malwarebytes.com/threat-analysis/2020/06/higaisa/**\n\nThreat Intelligence Team June 4, 2020\n\n_This post was authored by Hossein Jazi and Jérôme Segura_\n\nOn May 29th, we identified an attack that we believe is part of a new campaign from an\nAdvanced Persistent Threat actor known as Higaisa. The Higaisa APT is believed to be tied\nto the Korean peninsula, and was first disclosed by Tencent Security Threat Intelligence\nCenter in early 2019.\n\nThe group’s activities go back to at least 2016 and include the use of Trojans such as Gh0st\nand PlugX, as well as mobile malware. Its targets include government officials and human\nrights organizations, as well as other entities related to North Korea.\n\nIn this latest incident, Higaisa used a malicious shortcut file ultimately responsible for\ncreating a multi-stage attack that consists of several malicious scripts, payloads and decoy\nPDF documents.\n\n### Distribution\n\nThe threat actors used a malicious LNK file bundled within an archive file which was most\nlikely distributed via spear-phishing.\n\nWe were able to identify two variants of this campaign that possibly have been distributed\nbetween May 12th and 31st:\n\n\n-----\n\nCV_Colliers.rar\n“Project link and New copyright policy.rar”\n\nBoth RAR archives bundle two malicious LNK files. In the newer variant (CV_Colliers.rar),\nthe LNK files are disguised as a Curriculum Vitae (CV) and International English Language\nTesting System (IELTS) exam results. The older one (Project link and New copyright\n[policy.rar) seems to target product teams that are using zeplin.io.](https://zeplin.io/)\n\nThe following shows the overall process flow when executing the malicious LNK file.\n\nFigure 1: Process graph\n\n### LNK file\n\nThe LNK file contains a list of commands that will be executed upon running, and a blob that\nis a base64 encoded compressed payload. Here is the list of commands that will be\nexecuted:\n\nFigure 2: Malicious lnk commands\n\nCopy content of the LNK file into “g4ZokyumB2DC4.tmp” in %APPDATA% temp\ndirectory.\nCopy content of “certutil.exe” into “gosia.exe” ( “*ertu*.exe is used to bypass security\ndetection).\nLook for the base64 blob using “findstr.exe” and write it to “cSi1rouy4.tmp”.\n\n\n-----\n\nDecode content of cSi1rouy4.tmp using gosia.exe -decode (certutil.exe -decode)\nand write it to “o423DFDS4.tmp”.\nDecompress content of “o423DFDS4.tmp” in temp directory along with a decoy PDF\ndocument using “expand.exe -F:*” (Figure 3) .\nCopy “66DF3DFG.tmp” and “34fDKfSD38.js” files into “C:\\Users\\Public\\Downloads”\ndirectory.\nExecute the JS file by calling Wscript.\nOpen the decoy document.\n\nFigure 3: Content of the “o423DFDS4.tmp” cab file\nThe list of commands executed by this LNK shortcut is the same as the one reported by\n[Anomali on the Higasia Covid-19 campaign. The only difference is the name of the tmp files](https://www.anomali.com/blog/covid-19-themes-are-being-utilized-by-threat-actors-of-varying-sophistication)\nand name of certutil.exe which in this new case is “gosia.exe”, while in the March campaign\nthe name was “mosia.exe”.\n\nBoth LNK files embedded within the archive are executing similar commands with the\ndifferent Command and Control (C&C) configurations. Running each of them would show a\ndifferent decoy document.\n\n\n-----\n\nFigure 4: CV Decoy document\n\nFigure 5: IELTS test result decoy document\n\n### JS file\n\nThe JavaScript file performs the following commands:\n\nCreate “d3reEW.exe” in “C:\\Users\\Public\\Downloads” and store “cmd /c ipconfig” in it.\nExecute the dropped “svchast.exe”.\nCopy “svchhast.exe” into startup directory and rename it as “officeupdate.exe”.\n\n\n-----\n\nAdd officeupdate.exe to scheduled tasks.\nSend a POST request to a hardcoded URL with “d3reEW.exe” as data.\n\nFigure 6: JS content\n\nFigure7: POST request\n\n## svchast.exe\n\nSvchast.exe is a small loader that loads the content of the shellcode stored in\n“66DF3DFG.tmp”.\n\n\n-----\n\nFigure 8: Main function of svchast.exe\nIn fact, this shellcode is a wrapper around the final shellcode. It performs some checks and\nthen calls the final shellcode.\n\n\n-----\n\nFigure 9: Calling\n\nfinal shellcode\nThe final shellcode dynamically resolves the imports and allocates memory for the content\nthat will be executed.\n\nFigure 10: Resolving the imports\n\n\n-----\n\nFigure 11: Allocate memory for new thread\nFinally it calls “CreateThread” to create a thread within its memory space to make HTTPS\nrequests to its C&C server.\n\nFigure 11: CreateThread\nAt the time of analysis, the server was down so we weren’t able to clearly identify the\nultimate goal of this attack.\n\n### Chaining techniques for evasion\n\nWhile most malware campaigns use a simple decoy document that typically retrieves a\nmalware payload, more advanced attackers will often try unconventional means to infect\ntheir victims.\n\nWe reproduced this attack in our lab using an email as the infection vector, as we surmise\n[that victims were spear-phished. Malwarebytes (in this case the Nebula business version)](https://www.malwarebytes.com/business/cloud/)\nstopped the LNK file execution from WinRAR and therefore completely stopped the attack.\n\n\n-----\n\n### IOCs\n\n**CV_Colliers.rar**\n\ndf999d24bde96decdbb65287ca0986db98f73b4ed477e18c3ef100064bceba6d\n\n**Project link and New copyright policy.rar**\n\nc3a45aaf6ba9f2a53d26a96406b6c34a56f364abe1dd54d55461b9cc5b9d9a04\n\n**Curriculum Vitae_WANG LEI_Hong Kong Polytechnic University.pdf.lnk**\n\n50d081e526beeb61dc6180f809d6230e7cc56d9a2562dd0f7e01f7c6e73388d9\n\n**Tokbox icon – Odds and Ends – iOS – Zeplin.lnk**\n\n1074654a3f3df73f6e0fd0ad81597c662b75c273c92dc75c5a6bea81f093ef81\n\n**International English Language Testing System certificate.pdf.lnk**\n\nc613487a5fc65b3b4ca855980e33dd327b3f37a61ce0809518ba98b454ebf68b\n\n**Curriculum Vitae_WANG LEI_Hong Kong Polytechnic University.pdf.lnk**\n\ndcd2531aa89a99f009a740eab43d2aa2b8c1ed7c8d7e755405039f3a235e23a6\n\n\n-----\n\n**Conversations – iOS – Swipe Icons – Zeplin.lnk**\nc0a0266f6df7f1235aeb4aad554e505320560967248c9c5cce7409fc77b56bd5\n\n**C2 domains (ipconfig exfiltration)**\nsixindent[.]epizy[.]com\ngoodhk[.]azurewebsites[.]net\nzeplin[.]atwebpages[.]com\n\n**C2s used by svchast.exe**\n45.76.6[.]149\nwww.comcleanner[.]info\n\n### MITRE ATT&CK techniques\n\n**Tactic** **ID** **Name** **Details**\n\n\nExecution [T1059](https://attack.mitre.org/techniques/T1059) Command-Line\nInterface\n\n\nStarts CMD.EXE for commands\n(WinRAR.exe, wscript.exe) execution\n\n\n[T1106](https://attack.mitre.org/techniques/T1106) Execution through API Application (AcroRd32.exe) launched\nitself\n\n[T1053](https://attack.mitre.org/techniques/T1053) Scheduled Task Loads the Task Scheduler DLL interface\n(Officeupdate.exe)\n\n[T1064](https://attack.mitre.org/techniques/T1064) Scripting Executes scripts (34fDFkfSD38.js)\n\n[T1204](https://attack.mitre.org/techniques/T1204) User Execution Manual execution by user (opening LNK\nfile)\n\n\nPersistence [T1060](https://attack.mitre.org/techniques/T1060) Registry Run Keys /\nStartup Folder\n\n\nWrites to a start menu file\n(Officeupdate.exe)\n\n\nPrivilege\nEscalation\n\nDefense\nEvasion\n\n\n[T1053](https://attack.mitre.org/techniques/T1053) Scheduled Task Uses Task Scheduler to run other\napplications (Officeupdate.exe)\n\n[T1053](https://attack.mitre.org/techniques/T1053) Scheduled Task Uses Task Scheduler to run other\napplications (Officeupdate.exe)\n\n[T1064](https://attack.mitre.org/techniques/T1064) Scripting Executes scripts (34fDFkfSD38.js)\n\n\n[T1140](https://attack.mitre.org/techniques/T1140/) Deobfuscate/Decode\nFiles or Information\n\n\ncertutil to decode Base64 binaries,\nexpand.exe to decompress a CAB file\n\n\nDiscovery [T1012](https://attack.mitre.org/techniques/T1012) Query Registry Reads the machine GUID from the\nregistry\n\n\n-----\n\n[T1082](https://attack.mitre.org/techniques/T1082) System Information\nDiscovery\n\n[T1016](https://attack.mitre.org/techniques/T1016) System Network\nConfiguration\nDiscovery\n\n\nReads the machine GUID from the\nregistry\n\nUses IPCONFIG.EXE to discover IP\naddress\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2020/2020-06-03 - New LNK attack tied to Higaisa APT discovered.pdf"
    ],
    "report_names": [
        "2020-06-03 - New LNK attack tied to Higaisa APT discovered.pdf"
    ],
    "threat_actors": [
        {
            "id": "873919c0-bc6a-4c19-b18d-c107e4aa3d20",
            "created_at": "2023-01-06T13:46:39.138138Z",
            "updated_at": "2025-03-27T02:00:03.005534Z",
            "deleted_at": null,
            "main_name": "Higaisa",
            "aliases": [],
            "source_name": "MISPGALAXY:Higaisa",
            "tools": [],
            "source_id": "MISPGALAXY",
            "reports": null
        },
        {
            "id": "30c9c492-afc6-4aa1-8fe6-cecffed946e0",
            "created_at": "2022-10-25T15:50:23.400822Z",
            "updated_at": "2025-03-27T02:00:55.461334Z",
            "deleted_at": null,
            "main_name": "Higaisa",
            "aliases": [
                "Higaisa"
            ],
            "source_name": "MITRE:Higaisa",
            "tools": [
                "PlugX",
                "certutil",
                "gh0st RAT"
            ],
            "source_id": "MITRE",
            "reports": null
        },
        {
            "id": "c0cedde3-5a9b-430f-9b77-e6568307205e",
            "created_at": "2022-10-25T16:07:23.528994Z",
            "updated_at": "2025-03-27T02:02:09.847683Z",
            "deleted_at": null,
            "main_name": "DarkHotel",
            "aliases": [
                "APT-C-06",
                "ATK 52",
                "CTG-1948",
                "Dubnium",
                "Fallout Team",
                "Higaisa",
                "Luder",
                "Operation DarkHotel",
                "Operation Daybreak",
                "Operation Inexsmar",
                "Operation PowerFall",
                "Operation The Gh0st Remains the Same",
                "SIG25",
                "Shadow Crane",
                "T-APT-02",
                "Tungsten Bridge",
                "Zigzag Hail"
            ],
            "source_name": "ETDA:DarkHotel",
            "tools": [
                "Asruex",
                "DarkHotel",
                "DmaUp3.exe",
                "GreezeBackdoor",
                "Karba",
                "Nemain",
                "Nemim",
                "Ramsay",
                "Retro",
                "Tapaoux",
                "Trojan.Win32.Karba.e",
                "Virus.Win32.Pioneer.dx",
                "igfxext.exe",
                "msieckc.exe"
            ],
            "source_id": "ETDA",
            "reports": null
        }
    ],
    "ts_created_at": 1673536221,
    "ts_updated_at": 1743041741,
    "ts_creation_date": 1653708746,
    "ts_modification_date": 1653708746,
    "files": {
        "pdf": "https://archive.orkl.eu/e0974628ddea51d61fa5c71f4b880039f3b71cb7.pdf",
        "text": "https://archive.orkl.eu/e0974628ddea51d61fa5c71f4b880039f3b71cb7.txt",
        "img": "https://archive.orkl.eu/e0974628ddea51d61fa5c71f4b880039f3b71cb7.jpg"
    }
}