# 0-raw-data/ORKL/library/jsons_entries/40e2d4c1-5807-415c-8793-28b985a3df61.json
{
    "id": "40e2d4c1-5807-415c-8793-28b985a3df61",
    "created_at": "2023-04-05T02:08:59.808091Z",
    "updated_at": "2025-03-27T02:09:29.72205Z",
    "deleted_at": null,
    "sha1_hash": "4729aab25c0c9fb1b6f968cc96e970cc762e3eee",
    "title": "2023-02-28 - Can You See It Now- An Emerging LockBit Campaign",
    "authors": "",
    "file_creation_date": "2023-04-03T08:22:19Z",
    "file_modification_date": "2023-04-03T08:22:19Z",
    "file_size": 662094,
    "plain_text": "# Can You See It Now? An Emerging LockBit Campaign\n\n**[fortinet.com/blog/threat-research/emerging-lockbit-campaign](https://www.fortinet.com/blog/threat-research/emerging-lockbit-campaign)**\n\nFebruary 28, 2023\n\nFortiGuard Labs has observed a new LockBit ransomware campaign during last December and\nJanuary using a combination of techniques effective against AV and EDR solutions. LockBit has\nbeen one of the more dangerous ransomware, active since 2019. It was part of several\nsuccessful attacks against a large variety of industries, including critical infrastructure.\n\nThis blog post discusses the infection chain and Tactics, Techniques, and Procedures (TTPs) of\nthis campaign.\n\n## Overview\n\nDescriptions of the attack refer to the stages outlined in Figure 1 below. The attack starts with a\n.img container (1) and a social engineering technique of displaying a single file once it’s mounted\nwhile hiding the rest of its files from the user. It can also cause malware analysts to miss the\npayloads while examining the samples manually.\n\n\n-----\n\nFigure 1: Campaign execution stages\n\nThe user is then prompted to open the single visible shortcut (2) file.\n\nFIgure 2: Contents of the .img file, including the hidden files\n\nIn some of the cases that we’ve observed, a python script is executed (2.1) using the official\nPython embed package. The only purpose of the script is to run the subsequent BAT scripts.\n[Some variants used a known UAC bypass method abusing the legitimate fodhelper.exe (3.1).](https://v3ded.github.io/redteam/utilizing-programmatic-identifiers-progids-for-uac-bypasses)\nThis enables the attacker’s BAT file to run in a new elevated process without the user’s approval.\n\nFigure 3: UAC bypass implementation.\n\nThe BAT script (4) does several things:\n\n1.   Changes the password of the logged-in user.\n\n2.   Copies its files to C:\\ProgramData.\n\n3.   Ensures that after the system reboots, it logs in without user interaction (using SysInternals\nAutologon).\n\n4.   Tries to \na.   Set the next reboot to be in Safe Mode using bcdedit.exe.\n\nb.   Register a new service that will run its VBS script (4.1) using sc.exe.\n\nc.   Sets the service to run also in Safe-Mode using reg.exe.\n\n5.   If it fails, it sets in the registry a BAT file (4.2) to be run on logon as another UI shell by\nWinlogon.\n\n6.   Reboots the machine.\n\nFigure 4: Persistence of the BAT file (4.2).\n\nThe ransomware executable resides within a password-protected archive. The script that runs\nafter boot executes another BAT script (4.3) to extract the ransomware payload. It uses the 7-zip\narchiver and then runs it with a ‘-pass’ argument that is needed for the malicious executable to\nunpack itself.\n\nFigure 5: Command inside BAT script (4.3) decrypting and running the ransomware\n\n[The final payload is LockBit. Analysts from TrendMicro have published an analysis of the](https://www.trendmicro.com/en_us/research/22/g/lockbit-ransomware-group-augments-its-latest-variant--lockbit-3-.html)\nransomware.\n\nTargeting focuses on Spanish-speaking victims – all samples target Mexican or Spanish firms,\nmainly in the consulting and law sectors.\n\n\n-----\n\n**Ransomware note:**\n\n_Your data are stolen and encrypted_\n\n_The data will be published on TOR website if you do not pay the ransom_\n\n_You can contact us and decrypt one file for free on this TOR site_\n\n_(you should download and install TOR browser first https://torproject.org)_\n\n_http://<onion address>_\n\n_Your company id for log in: <unique id>_\n\n## Evasive Tradecraft\n\nThe detection rate of the samples in VirusTotal was a minimal single digit, with some completely\nundetected, suggesting the campaign’s methods are effective in defense evasion.\n\nFigure 6: Detection in VirusTotal of one of the .img files.\n\nDelivery through a .img container bypasses the Mark of The Web (MOTW) protection\nmechanism. Multi-stage scripts that extract a password-protected ransomware executable, which\nis unpacked only when run with a unique password, allow evading traditional signature-based\ndetection.\n\nThe malware authors have shown a creative and wide-ranging usage of signed, legitimate\nexecutables: the mounting of .img files by Windows Explorer, python execution by a signed\ninterpreter, the extraction of encrypted archives by 7-zip, and automatic log-in using Sysinternals’\nAutologon. This allows for minimal reliance on custom code, trimming development costs, and\nstaying under the radar of EDRs.\n\n## Summary\n\nThis campaign’s highly evasive nature demonstrates that attackers continue to leverage\nincreasingly obscure methodologies to avoid detection. This unique combination of executables\nwasn’t seen in previous LockBit attacks.\n\n[These payloads may be related to the LockBit builder leak in late 2022. Hence, a definitive](https://www.malwarebytes.com/blog/news/2022/09/lockbit-builder-leaked-by-disgruntled-developer)\nattribution to the original group behind LockBit, or its affiliates, may be difficult. Other\ncybercriminals would want to associate themselves with the deterrence already afforded by past\nacts of LockBit.\n\n### Fortinet Solutions\n\n\n-----\n\nFortiEDR detects and blocks these threats out of the box without any prior knowledge or special\nconfiguration. It does this using its post-execution prevention engine to identify malicious\nactivities:\n\nFigure 7: FortiEDR blocking the ransomware.\n\nAll network IOCs have been added to the FortiGuard WebFilter blocklist.\n\nFortiGuard Antivirus has coverage in place as follows:\n\nW32/Lockbit.K!tr.ransom\n\nThe FortiGuard Antivirus service engine is included in Fortinet’s FortiGate, FortiMail, FortiClient,\nand FortiEDR solutions.\n\n[In addition, as part of our membership in the Cyber Threat Alliance, details of this threat were](https://www.fortiguard.com/cta)\nshared in real-time with other Alliance members to help create better protections for customers.\n\nLearn more about Fortinet’s [FortiGuard Labs threat research and intelligence organization and](https://www.fortinet.com/fortiguard/labs?utm_source=blog&utm_medium=blog&utm_campaign=fglabs)\n[the FortiGuard AI-powered security services portfolio. If you think this or any other cybersecurity](https://www.fortinet.com/solutions/enterprise-midsize-business/security-as-a-service/fortiguard-subscriptions?utm_source=blog&utm_medium=blog&utm_campaign=services)\n[threat has impacted your organization, contact our Global FortiGuard Incident Response Team.](https://www.fortinet.com/corporate/about-us/contact-us/experienced-a-breach)\n\n## Appendix A: MITRE ATT&CK Tactics and Techniques\n\n**Tactic \\ Technique ID** **Description**\n\n**TA0002** **Execution**\n\nT1059.003 Command and Scripting Interpreter: Windows Command Shell\n\nT1059.005 Command and Scripting Interpreter: Visual Basic\n\nT1059.006 Command and Scripting Interpreter: Python\n\nT1204.002 User Execution: Malicious File\n\n**TA0003** **Persistence**\n\nT1547.004 Boot or Logon Autostart Execution: Winlogon Helper DLL\n\nT1053.005 Scheduled Task/Job: Scheduled Task\n\n\n-----\n\nT1543.003 Create or Modify System Process: Windows Service\n\n**TA0005** **Defense Evasion**\n\nT1553.005 Subvert Trust Controls: Mark-of-the-Web Bypass\n\nT1548.002 Abuse Elevation Control Mechanism: Bypass User Account Control\n\nT1027.002 Obfuscated Files or Information: Software Packing\n\nT1562.009 Impair Defenses: Safe Mode Boot\n\n**TA0040** **Impact**\n\nT1486 Data Encrypted for Impact\n\nT1529 System Shutdown/Reboot\n\n## APPENDIX B: IoCs\n\n**IMG SHA256**\n\n1ef3ae251833be08b6f3e525969ae02c28cb0238e3adb3091e572a10633f7ef7\n\ndad61d9f919a9cc84ae633e948946e7546b21dc4d9d47d19d96fd308c7de40cb\n\nd73bcd2e29191b260a26d87c3035bde33163cc319649291db9f04c48c94da896\n\nee2b182a56ded459a113513985ff624631a9515c7efa2282708483ace640eb3a\n\ndca325a0028dc8e41dcf739cd00701a19066fc88c0d22be5316f7a4b7b219fe8\n\n35bf036bf46fa21f3354d60a2c50d2959e1e9193bec8364575dc3fd4644732ae\n\n781ead305cdb5fa0153369431dedd40fe138308fcdf5dfda1cfeaaba296752e3\n\n1ef3ae251833be08b6f3e525969ae02c28cb0238e3adb3091e572a10633f7ef7\n\ndad61d9f919a9cc84ae633e948946e7546b21dc4d9d47d19d96fd308c7de40cb\n\nd73bcd2e29191b260a26d87c3035bde33163cc319649291db9f04c48c94da896\n\nee2b182a56ded459a113513985ff624631a9515c7efa2282708483ace640eb3a\n\ndca325a0028dc8e41dcf739cd00701a19066fc88c0d22be5316f7a4b7b219fe8\n\n35bf036bf46fa21f3354d60a2c50d2959e1e9193bec8364575dc3fd4644732ae\n\n781ead305cdb5fa0153369431dedd40fe138308fcdf5dfda1cfeaaba296752e3\n\n1858a862390adcaa4cea6782e7dba077697475ff9ada9d75c4897ccd563998af\n\n**Ransomware Executables**\n\n\n-----\n\nSHA256 Name\n\ncb049c6e59106bbdfd804a9d02bb31ea09a3918018cbb97fb12d2bcf9e475465 documentos.exe\n\n334148a7434f4fd27dcc6600edc2f29e4f11ada0be9f71f807cbd4154abd74be documentos.exe\n\nfd3577ff36496320485ffa05681ffa516a56fc4818c3fc89774aa4bb20e2c17f documentos.exe\n\n8465c979990e75262d15e93453287d6107f008035d6d6a05bd3a92c2e3fe1d40 HacAK.exe\n\n40828437094a02ab467a0c0343d08c110d3b0c2972b609bcdd355667614209f EMJgp.exe\n\n50f49ac742a127085e0a824bcae7e25326ea0ef0741f0abe34ce494f2c4ef4d2 byhHI.exe\n\ncc58dcd32a440e7f95d19b653a55c1e2c383efc2bd19443238dd3008c1cbe147 bOpDX.exe\n\n6eb6431dcfb1e7105fb05e2d8b01e231f6d4b82a1df3639499d0adacd00757cc gVozH.exe\n\n**Domains**\n\npoliovocalist[.]com\n\n**IPs**\n\n198.244.187[.]248\n\n150.129.218[.]231\n\n**LockBit Portal URLs**\n\nhxxp://lockbit3jx6je7tm6hhm6zzafgy6hpil3ur6jmc2a4ugan7xzztv6oqd[.]onion\n\nhxxp://lockbitdvbpfczc3yrs37kpp6avnrgr7yygi2f45qxvef2yqi36lpxyd[.]onion\n\nhxxp://lockbit3hc6syym13ki2ag5jskr6q5qa3spspjpmtfhh6fufut737zid[.]onion\n\nhxxp://lockbitov3afmxgknfhk2o5d4uqrhygd7ty3xqm56qd6zjlu6u43pgyd[.]onion\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2023/2023-02-28 - Can You See It Now- An Emerging LockBit Campaign.pdf"
    ],
    "report_names": [
        "2023-02-28 - Can You See It Now- An Emerging LockBit Campaign.pdf"
    ],
    "threat_actors": [
        {
            "id": "aa73cd6a-868c-4ae4-a5b2-7cb2c5ad1e9d",
            "created_at": "2022-10-25T16:07:24.139848Z",
            "updated_at": "2025-03-27T02:02:10.120505Z",
            "deleted_at": null,
            "main_name": "Safe",
            "aliases": [],
            "source_name": "ETDA:Safe",
            "tools": [
                "DebugView",
                "LZ77",
                "OpenDoc",
                "SafeDisk",
                "TypeConfig",
                "UPXShell",
                "UsbDoc",
                "UsbExe"
            ],
            "source_id": "ETDA",
            "reports": null
        }
    ],
    "ts_created_at": 1680660539,
    "ts_updated_at": 1743041369,
    "ts_creation_date": 1680510139,
    "ts_modification_date": 1680510139,
    "files": {
        "pdf": "https://archive.orkl.eu/4729aab25c0c9fb1b6f968cc96e970cc762e3eee.pdf",
        "text": "https://archive.orkl.eu/4729aab25c0c9fb1b6f968cc96e970cc762e3eee.txt",
        "img": "https://archive.orkl.eu/4729aab25c0c9fb1b6f968cc96e970cc762e3eee.jpg"
    }
}