# 0-raw-data/ORKL/library/jsons_entries/45576bf8-3c3c-498d-b4df-20e482ce4559.json
{
    "id": "45576bf8-3c3c-498d-b4df-20e482ce4559",
    "created_at": "2023-05-06T02:08:57.526503Z",
    "updated_at": "2025-03-27T02:09:29.447361Z",
    "deleted_at": null,
    "sha1_hash": "8ee10a9efb13f0b7bc62e93c7991a818111f1a62",
    "title": "2023-03-24 - Aurora- The Dark Dawn and its Menacing Effects",
    "authors": "",
    "file_creation_date": "2023-05-05T01:53:21Z",
    "file_modification_date": "2023-05-05T01:53:21Z",
    "file_size": 960560,
    "plain_text": "# Aurora: The Dark Dawn and its Menacing Effects\n\n**[research.loginsoft.com/threat-research/aurora-the-dark-dawn-and-its-menacing-effects/](https://research.loginsoft.com/threat-research/aurora-the-dark-dawn-and-its-menacing-effects/)**\n\n**Aurora Stealer Background**\n\nMarch 24, 2023\n\nBy Saharsh Agrawal\n\n\nMarch 23, 2023\n\n\n[Aurora Stealer a non-friendly program in your neighborhood is developed in Go-lang. The](https://go.dev/)\nmalware was developed for threat actors to steal victims’ sensitive data. Since its initial\nrelease, there have been many adaptations to the malware. It started as a botnet having\ninfo-stealer capabilities but currently it is a full-time info-stealer which is being sold by the\nthreat actor labelled Cheshire on Russian speaking dark-net forums mentioned by\n[cybersecurity firm Sekoia.](https://blog.sekoia.io/aurora-a-rising-stealer-flying-under-the-radar/)\n\nThe malware has a knack for stealing the victim’s browser data including cookies,\npasswords, login data and many more. It has the horizon to steal data from most browsers.\nThe stealer also targets crypto wallets to pilfer victim’s data. The malware does not stop\nhere, it can exfiltrate files and data stored in different software from the victim’s machine.\nOnce all the super-sensitive data is possessed by the stealer it delivers the data to the\nattacker.\n\n**Operating Model**\n\n_Figure1: Attack_ _Chain of Aurora Stealer_\n\n\n-----\n\nThe malware is being spread among victims in the hospitality industry through phishing\n[emails as seen by Trendmirco. Besides phishing campaigns, the threat actors also have](https://www.trendmicro.com/vinfo/us/threat-encyclopedia/spam/3721/fake-roadway-map-leads-to-aurora-stealer)\nutilized SEO poisoning technique to deliver the malware via malicious ads of notepad++\n[according to SANS finding.](https://isc.sans.edu/diary/Malicious+Google+Ad+Fake+Notepad+Page+Aurora+Stealer+malware/29448)\n\nUpon the detonation of the malware, our observations revealed the implementation of\nmultiple sophisticated techniques by the malware to evade detection by antivirus software.\n[One such technique involved the incorporation of binary padding, resulting in an increase in](https://attack.mitre.org/techniques/T1027/001/)\nthe file size allowing the malware to circumvent the antivirus software and evade detection.\n\n\n-----\n\n_Figure2: Fake Adobe Acrobat Installer Downloaded from SEO Poisoning Campaign_\n\nAfter successful infiltration, the malware utilizes Windows Management Instrumentation\n(WMI) commands to gather pertinent system information. This behavior was consistently\nobserved in all samples of Aurora stealer that we analyzed. Additionally, independent\n[research conducted by Cyble and other industry experts has reported this observation.](https://blog.cyble.com/2023/01/18/aurora-a-stealer-using-shapeshifting-tactics/)\n\n_Figure2.1 : WMI Commands getting Executed_\nWe also captured Aurora stealer’s actions that were aimed at extracting sensitive information\nfrom files and folders located in various browsers, cryptocurrency wallets, and drives, as\nillustrated below.\n\n\n-----\n\n_Figure3: Stealing Browser and Crypto Data_\nAurora leverages PowerShell to facilitate the execution of its second stage, employing the\nfollowing command:\n```\n powershell \"\" \"start-process C:\\Users\\Admin\\AppData\\Local\\Temp\\BK2c7L9Si8.exe\"\n\n```\nDuring the second stage of Aurora, an entry to the “run keys” in the registry was observed,\nenabling the malware to automatically execute its third stage DLL via RUNDLL32.exe at\neach user login.\n\n_Figure4: Run Key Added for DLL Execution_\nThe third stage of the malware executes an encoded PowerShell command allowing it to\ndelay execution for certain duration using the following command before it exfiltrates the data\nto the attacker.\n```\nC:\\Windows\\System32\\WindowsPowerShell\\v1.0\\powershell.exe\" -enc\nUwB0AGEAcgB0AC0AUwBsAGUAZQBwACAALQBTAGUAYwBvAG4AZABzACAAMwA0AA==\n\n```\nOnce all the data and information are collected, the malware employs advanced encoding\ntechniques to conceal the information and prepare it for shipment to the attacker over the C2\nchannel. The attacker can retrieve the encoded data and use it for nefarious purposes, such\nas identity theft or financial fraud.\n\n\n-----\n\n_[Figure5: C2 Panel; source: twitter](https://twitter.com/silascutler/status/1614518491557642243)_\n\n**Detection**\n\nThe infernal behavior of this malware can be detected using SIGMA rules.\n\nTo begin with, we attribute different WMI commands getting executed in sequence to detect\nAurora Stealer malware.\n```\nlogsource:\n\n category: process_creation\n\n product: windows\n\ndetection:\n\n selection:\n\n Image|endswith: '\\wmic.exe'\n\n CommandLine|contains:\n\n - 'wmic os get Caption'\n\n - 'wmic path win32_VideoController get name'\n\n - 'wmic cpu get name'\n\n condition: selection\n\n\n```\nOne detection would be to capture changes in ‘run key’ registry which execute DLL using\nRUNDLL32.exe\n\n\n-----\n\n```\nlogsource:\n\n product: windows\n\n category: registry_set\n\ndetection:\n\n selection:\n\n EventType: SetValue\n\n TargetObject|contains:\n'\\REGISTRY\\MACHINE\\Software\\Microsoft\\Windows\\CurrentVersion\\RunOnce'\n\n Details|re: .*rundll32.exe.C:\\\\\\\\Windows\\\\\\\\system32\\\\\\\\advpack.dll,DelNodeRunDLL32\n\\\\\"C:\\\\\\\\Users\\\\\\\\.*\\\\AppData\\\\\\\\Local\\\\\\\\Temp\\\\\\\\IXP[0-9]{1,3}.TMP.*\n\n condition: selection\n\n```\nBased on the general behavior followed by Aurora and different malwares, we have\ndeveloped several detection methods. We also observed changes in the configuration of\nWindows Defender, where the path ‘C:\\Program Data’ is excluded.\n```\nlogsource:\n\n product: windows\n\n category: registry_set\n\ndetection:\n\n selection:\n\n EventType: SetValue\n\n TargetObject|contains: '\\Microsoft\\Windows Defender\\Exclusions'\n\n Details: 'C:\\Program Data\\'\n\n condition: selection\n\n\n```\nAnother detection would be by tracing event id `5007`, which can be activated by enabling\nWindows Defender operation logs. This particular event is triggered whenever there is a\nmodification made to the configuration of Microsoft Defender Antivirus.\n```\nlogsource:\n\n product: windows\n\n service: windefend\n\ndetection:\n\n selection:\n\n EventID: 5007\n\n NewValue|contains: 'C:\\Program Data\\'\n\n condition: selection\n\n\n```\nAdditionally, one clever way to detect this malware is to monitor the execution of obfuscated\nPowerShell commands.\n\n[The same campaign can be detected using OSQUERY rules. Below are some rules that will](https://www.loginsoft.com/blog/2020/06/25/explicating-the-concepts-of-osquery/)\nhelp you keep your organization safe.\n\n\n-----\n\n```\nquery_wmic:\n\n  SELECT name,\n\n      pid,\n\n      cmdline,\n\n      path,\n\n      parent\n\n  FROM processes\n\n  WHERE LOWER(name) = 'wmic.exe'\n\n  AND\n\n  (\n\n   cmdline LIKE '%OS Get Caption%'\n\n   OR cmdline LIKE '%path win32_VideoController get name%'\n\n   OR cmdline LIKE '%cpu get name%'\n\n  );\n\n query_regsitry:\n\n  SELECT key,\n\n      path,\n\n      name,\n\n      type,\n\n      data\n\n  FROM registry\n\n  WHERE path LIKE\n'HKEY_LOCAL_MACHINE\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunOnce\\\\%'\n\n  AND data LIKE '%rundll32.exe C:\\\\Windows\\\\system32\\\\advpack.dll,DelNodeRunDLL32\n\"C:\\\\Users\\\\%\\\\AppData\\\\Local\\\\Temp\\\\IXP%.TMP\\\\%';\n\nquery_exclusion:\n\n  SELECT key,\n\n      path,\n\n      name,\n\n      type,\n\n      data\n\n  FROM registry\n\n  WHERE path LIKE 'HKEY_LOCAL_MACHINE\\SOFTWARE\\Microsoft\\Windows\nDefender\\Exclusions\\Paths\\%'\n\n  AND\n\n  (\n\n   data LIKE '%c:\\programdata%'\n\n   OR data LIKE '%\\appdata\\roaming\\microsoft\\windows\\start menu\\programs\\startup%'\n\n```\n[For more threat analytics reach us here.](https://research.loginsoft.com/contact-us/)\n\n[The above rules have been vindicated by simulating Aurora stealer in a sandbox](https://tria.ge/221117-p1yg5aad8t/behavioral1)\nenvironment.\n\n**MITRE ATT&CK Techniques**\n\n|Tactic|Technique ID|Technique Name|\n|---|---|---|\n\n\n-----\n\n|Execution|T1204  T1059.001  T1059.003  T1047|User Execution  Command and Scripting Interpreter: PowerShell  Command and Scripting Interpreter: Windows Command Shell  Windows Management Instrumentation|\n|---|---|---|\n|Persistence|T1547.001|Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder|\n|Defense Evasion|T1562.001 T1027 T1218 .011 T1497|Impair Defenses Obfuscated Files or Information System Binary Proxy Execution: R undll32 Virtualization/Sandbox Evasion|\n|Credential Access|T1555.003 T1056 T1003 T1528 T1539 T1552 .002|Credentials from Password Stores: Credentials from Web Browsers Input Capture OS Credential Dumping Steal Application Access Token Steal Web Session Cookie Credentials in Registry|\n|Discovery|T1087 T1083 T1012 T1082 T1518 .001 T1614 T1497|Account Discovery File and Directory D iscovery Query Registry System Informa tion Discovery Software Discovery: Security S oftware Discovery System Location Discovery Virtualization/Sandbox Evas ion|\n|Collection|T1119 T1005 T1056 T1113|Automated Collection Data from Local Syste m Input Capture Screen Captu re|\n|Command and Control|T1071.001 T1105 T1095 T1571|Application Layer Protocol: Web Protocols Ingress Tool Transfer Non-Application Laye r Protocol Non-Standard Port|\n\n\n**Threat Bites**\n\n**Threat Actor**\n\n**Targeted Country**\n\n**Targeted Industry**\n\n**First Seen**\n\n**Last Seen**\n\n\n-----\n\n**LOLBAS**\n\n**Telemetry**\n\n**Samples**\n\n:\n\n:\n\n:\n\n:\n\n:\n\n:\n\n:\n\n:\n\nUnknown\n\nWorldwide\n\nHotels\n\nApril 2022\n\nMarch 2023\n\nWmic, Reg, Rundll32\n\nSysmon, Security, Windefend, PowerShell\n\n[https://tria.ge/230118-f1ewcaac94/behavioral1](https://tria.ge/230118-f1ewcaac94/behavioral1)\n\n**Author: Saharsh Agrawal**\n\nSecurity Researcher, Loginsoft\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2023/2023-03-24 - Aurora- The Dark Dawn and its Menacing Effects.pdf"
    ],
    "report_names": [
        "2023-03-24 - Aurora- The Dark Dawn and its Menacing Effects.pdf"
    ],
    "threat_actors": [
        {
            "id": "aa73cd6a-868c-4ae4-a5b2-7cb2c5ad1e9d",
            "created_at": "2022-10-25T16:07:24.139848Z",
            "updated_at": "2025-03-27T02:02:10.120505Z",
            "deleted_at": null,
            "main_name": "Safe",
            "aliases": [],
            "source_name": "ETDA:Safe",
            "tools": [
                "DebugView",
                "LZ77",
                "OpenDoc",
                "SafeDisk",
                "TypeConfig",
                "UPXShell",
                "UsbDoc",
                "UsbExe"
            ],
            "source_id": "ETDA",
            "reports": null
        }
    ],
    "ts_created_at": 1683338937,
    "ts_updated_at": 1743041369,
    "ts_creation_date": 1683251601,
    "ts_modification_date": 1683251601,
    "files": {
        "pdf": "https://archive.orkl.eu/8ee10a9efb13f0b7bc62e93c7991a818111f1a62.pdf",
        "text": "https://archive.orkl.eu/8ee10a9efb13f0b7bc62e93c7991a818111f1a62.txt",
        "img": "https://archive.orkl.eu/8ee10a9efb13f0b7bc62e93c7991a818111f1a62.jpg"
    }
}