# 0-raw-data/ORKL/library/jsons_entries/46508a11-c105-4bcd-acf9-f7e257dbbb41.json
{
    "id": "46508a11-c105-4bcd-acf9-f7e257dbbb41",
    "created_at": "2023-01-12T15:02:42.74443Z",
    "updated_at": "2025-03-27T02:06:13.110348Z",
    "deleted_at": null,
    "sha1_hash": "e9dc2ddf2cfda0ac448a34b02460fef4531c6c28",
    "title": "2021-01-26 - Cybereason vs. RansomEXX Ransomware",
    "authors": "",
    "file_creation_date": "2022-05-28T00:47:04Z",
    "file_modification_date": "2022-05-28T00:47:04Z",
    "file_size": 760731,
    "plain_text": "# Cybereason vs. RansomEXX Ransomware\n\n**[cybereason.com/blog/cybereason-vs.-ransomexx-ransomware](https://www.cybereason.com/blog/cybereason-vs.-ransomexx-ransomware)**\n\nOver the last few months, the [Cybereason Nocturnus Team has been tracking the activity](https://www.cybereason.com/company/nocturnus)\n[around the RansomEXX ransomware. It has been active since 2018, but came to fame in](https://malpedia.caad.fkie.fraunhofer.de/details/win.ransomexx)\n2020 in attacks on major organizations such as the Texas Department of Transportation.\n[RansomEXX started as a Windows variant, but a Linux variant was discovered earlier this](https://securelist.com/ransomexx-trojan-attacks-linux-systems/99279/)\nyear.\n\n## Key Findings\n\n**Human-operated targeted attacks: RansomEXX is being used as a part of multi-staged**\nhuman-operated attacks targeting various government related entities and tech companies. It\nis being delivered as a secondary payload after initial compromise of the targeted network.\n\n**[Disables security products: The Windows variant has a functionality that was seen before](https://www.cybereason.com/blog/cybereason-vs.-clop-ransomware)**\nin other ransomware, disabling various security products for a smooth execution on the\ninfected machine.\n\n**Multi-Platform: RansomEXX started solely as a Windows variant, but later a Linux variant**\nwas added to the arsenal, sharing similarities with its predecessor.\n\n**Fileless ransomware: RansomEXX is usually delivered as a secondary in-memory payload**\nwithout ever touching the disk, which makes it harder to detect.\n\n\n-----\n\n**Detected and prevented: The** [Cybereason Defense Platform fully detects and prevents the](https://www.cybereason.com/platform#graphic)\nRansomEXX ransomware.\n\n## Background\n\nTheRansomEXX family, also known as Defray777 and Ransom X, runs as a solely inmemory payload that is not dropped to disk, making it highly evasive. RansomEXX was\n[involved in three major attacks in 2020 against Texas TxDOT in May of 2020, against Konica](https://www.bleepingcomputer.com/news/security/new-ransom-x-ransomware-used-in-texas-txdot-cyberattack/)\n[Minolta in the end of July, and against Brazil's court system in the beginning of November.](https://www.bleepingcomputer.com/news/security/brazils-court-system-under-massive-ransomexx-ransomware-attack/)\n\nIn addition, last December RansomEXX operators published stolen credentials from\nEmbraer, one of the largest aircraft makers in the world, on its own leaks website as part of\n[the ongoing double extortion trend.](https://www.infosecurity-magazine.com/blogs/double-extortion-ransomware/)\n\nIn mid 2020, a Linux variant of RansomEXX emerged. This variant, despite sharing\n[similarities with the Windows variant, is simpler than its predecessor and lacks many](https://securelist.com/ransomexx-trojan-attacks-linux-systems/99279/)\nfeatures such as disabling security software and command and control communication.\nThere are decryptors for both variants, and the threat actors send paying victims a private\nkey to decode their files.\n\n## RansomEXX Analysis\n\nThis analysis focuses on the Windows variant of RansomEXX, which can be classified as\n[fileless malware because it is reflectively loaded and executed in memory without touching](https://www.andreafortuna.org/2017/12/08/what-is-reflective-dll-injection-and-how-can-be-detected/)\nthe disk. Analysis of this sample reveals that it is partially obfuscated but includes indicative\ninformation such as the “ransome.exx” string that can be seen hard coded in the binary:\n\n_ransom.exx string hardcoded in the binary_\n\nUpon execution, RansomEXX starts decrypting some strings necessary for its operation:\n\n\n-----\n\n_RansomEXX’s strings decryption routine_\n\nThe mutex the malware creates is generated from the GUID of the infected machine:\n\n_The GUID generated on the infected machine_\n\nThe decrypted strings at this point include mainly logs:\n\n_Decrypted logging string_\n\nRansomEXX spawns a separate thread in the background to handle the logging process.\n\nWhen debugging the sample, the logs themselves can be seen in the console:\n\n_Logging as seen in the command line_\n\n[The malware then continues with terminating processes and system services that may](https://unit42.paloaltonetworks.com/vatet-pyxie-defray777/3/)\ninterfere with the execution, but excludes those that are relevant for its execution:\n\n\n-----\n\n_Some of the terminated services as well as processes excluded from termination_\n\nCybereason detects the execution of RansomEXX together with the below listed commands\nthat are executed post-encryption. These commands’ role is to prevent the victim from\nrestoring their system by deleting backups, Windows error recovery etc. Cybereason also\ndetects this malicious usage of Windows utilities:\n\n\n-----\n\n_RansomEXX’s attack tree as seen in the Cybereason Defense Platform_\n\nThe depicted above commands are as follows:\n\n**Command** **Action**\n\n\n\"C:\\Windows\\System32\\fsutil.exe\" usn deletejournal /D\nC:\n\n\"C:\\Windows\\System32\\wbadmin.exe\" delete catalog quiet\n\n\"C:\\Windows\\System32\\wevtutil.exe\" cl Setup\n\n\"C:\\Windows\\System32\\wevtutil.exe\" cl System\n\n\"C:\\Windows\\System32\\wevtutil.exe\" cl Application\n\n\"C:\\Windows\\System32\\wevtutil.exe\" cl Security\n\n\"C:\\Windows\\System32\\bcdedit.exe\" /set {default}\nbootstatuspolicy ignoreallfailures\n\n\"C:\\Windows\\System32\\bcdedit.exe\" /set {default}\nrecoveryenabled no\n\n\nfsutil.exe deletes the Update\nSequence Number journal\n\nwbadmin.exe deletes the\n[backup catalog](https://docs.microsoft.com/en-us/windows-server/administration/windows-commands/wbadmin-delete-catalog)\n\n[wevtutil clears event logs](https://winaero.com/how-to-clear-the-windows-event-log-from-the-command-line/)\n\n[bcdedit disable recovery mode](https://docs.microsoft.com/en-us/windows-hardware/drivers/devtest/bcdedit--set)\n\n\n-----\n\n\"C:\\Windows\\System32\\cipher.exe\" /w:C: cipher [overwrites deleted data](https://www.thewindowsclub.com/cipher-command-line-tool-windows)\nin drive C\n\n\n\"C:\\Windows\\System32\\schtasks.exe\" /Change /TN\n\"\\Microsoft\\Windows\\SystemRestore\\SR\" /disable\n\n\nschtasks disables the system\nrestore scheduled task\n\n\n\"C:\\Windows\\System32\\wevtutil.exe\" sl Security /e:false wevtutil disables the security\nevent [logs](https://vimalshekar.github.io/scriptsamples/Advanced-ETW-logging)\n\nAfter preparation of the environment RansomEXX encrypted the files on the victim’s machine\nand the following note is left on the machine:\n\n_The ransom note left on the victim’s machine_\n\nThe commands that disable file recovery and system restore after successfully encrypting\nthe victim’s files, and can also be observed clearly in the sample’s code:\n\n\n-----\n\n_Part of the post-encryption commands in RansomEXX’s code_\n\n## Cybereason Detection and Prevention\n\nCybereason detects the Windows utilities that are executed post-encryption as malicious and\n[triggers a Malop(™) for all of them:](https://www.cybereason.com/platform#malop)\n\n_Detection of the ransomware and malicious uses of windows utilities by the Cybereason_\n_Defense Platform_\n\nLooking at the Malop that was triggered by fsutil, the evidence for malicious activity can be\nseen together with the suspicions mapped to the MITRE ATT&CK matrix:\n\n\n-----\n\n_Suspicions and evidence triggered by fsutil_\n\nWhen Cybereason anti-ransomware prevention is turned on, the execution of the\nRansomEXX is prevented using the AI module:\n\n_Execution prevention of RansomEXX by the Cybereason Defense Platform_\n\n## Security Recommendations\n\n\n-----\n\n**• Enable the Anti-Ransomware Feature on Cybereason NGAV: Set Cybereason Anti-**\nRansomware protection mode to Prevent - more information for customers can be found\nhere\n\n**• Enable Anti-Malware Feature on Cybereason NGAV: Set Cybereason Anti-Malware**\nmode to Prevent and set the detection mode to Moderate and above - more information can\nbe found here\n\n**• Keep Systems Fully Patched: Make sure your systems are patched in order to mitigate**\nvulnerabilities\n\n**• Regularly Backup Files to a Remote Server: Restoring your files from a backup is the**\nfastest way to regain access to your data\n\n**• Use Security Solutions: Protect your environment using organizational firewalls, proxies,**\nweb filtering, and mail filtering\n\n**• Indicator's of Compromise: Includes C2 Domains, IP addresses, Docx files SHA-1**\nhashes, and Msi files. Open the chatbot on the lower right-hand side of this blog to download\nyour copy.\n\n## MITRE ATT&CK BREAKDOWN\n\n\n**Defense**\n**Evasion**\n\n**Impair**\n**Defenses:**\n**Disable or**\n**Modify Tools**\n\n**Indicator**\n**Removal on**\n**Host: File**\n**Deletion**\n\n\n**Impact** **Execution** **Discovery** **Privilege**\n**Escalation**\n\n\n**Data**\n**[Encrypted](https://attack.mitre.org/techniques/T1486/)**\n**for Impact**\n\n**Inhibit**\n**[System](https://attack.mitre.org/techniques/T1490/)**\n**Recovery**\n\n\n**Command and**\n**Scripting Interpreter:**\n**Windows Command**\n**Shell**\n\n**Command and**\n**[Scripting Interpreter:](https://attack.mitre.org/techniques/T1059/004/)**\n**Unix Shell**\n\n\n**Obfuscated Files**\n**or Information**\n\n**System**\n**[Information](https://attack.mitre.org/techniques/T1082/)**\n**Discovery**\n\n\n**Process**\n**Injection**\n\n\n**[Scheduled Task/Job](https://attack.mitre.org/techniques/T1053/)** **File and Directory**\n**Discovery**\n\n**Software**\n**Discovery:**\n**Security Software**\n**Discovery**\n\n\n-----\n\n**Process**\n**Discovery**\n\n**Daniel Frank**\n\nDaniel Frank is a senior Malware Researcher at Cybereason. Prior to Cybereason, Frank\nwas a Malware Researcher in F5 Networks and RSA Security. His core roles as a Malware\nResearcher include researching emerging threats, reverse-engineering malware and\ndeveloping security-driven code. Frank has a BSc degree in information systems.\n\n## RansomEXX Ransomware | Indicator's of Compromise\n\n**IOC** **Type** **Description**\n\n\n-----\n\n0abaa05da2a05977e0baf68838cff1712f1789e0\n\n6fae9aa52fd89bac83b69c2fbdc65c96e886427f\n\n06606fea0daaa99bd8ebfeb60f19976c20e6bb72\n\n0122efe580848879bb70f40ede63cb2edbfb4163\n\nccfc9578f721fbad30aa74facf20817abe118bfd\n\n423a2bf7ac322273bdacf638703ea99c44462862\n\ndfc37340f5deaa89681539b0f5c22059aac4c31d\n\n9711cdf002e5b7ecccfa309058d53dde67b029ee\n\n3e6689dc6a8a717b4114a7fe65bba594c597c7b9\n\n18b2704b49828035148aebe9e77b286a30c702b6\n\ne7748b92347f95589fa739cbe5c089046614ce92\n\n427178528152670c68f2f2937f05a5cdfebff1c2\n\n3555aaebe6c113fb8f923a38cb3bd75da6e86277\n\n6185e3514a32d2f3fb9ce292ba514d01584cced8\n\nfc9284b7a140c0d411ebd0eb4752e477d5d213fc\n\n11eec31710902820e79ba1e363d4c1256b75c615\n\n5238ba19bb3c7298ee13fe6eb0cf5f8787c13cd8\n\n24e773aa271fc0636cda6b0966a6034b65cb3052\n\n91ad089f5259845141dfb10145271553aa711a2b\n\n132def0d906a53360bdbdd3da109bfa41bcdbb6c\n\n3bf79cc3ed82edd6bfe1950b7612a20853e28b09\n\n50f191f04aa6cff1d8688a3c5d6cce96739ab6b3\n\n\nSHA1 RansomEXX Windows\nExecutable\n\nSHA1 RansomEXX Linux\nExecutable\n\n\n-----\n\nAbout the Author\n\n**Cybereason Nocturnus**\n\nThe Cybereason Nocturnus Team has brought the world’s brightest minds from the military,\ngovernment intelligence, and enterprise security to uncover emerging threats across the\nglobe. They specialize in analyzing new attack methodologies, reverse-engineering malware,\nand exposing unknown system vulnerabilities. The Cybereason Nocturnus Team was the first\nto release a vaccination for the 2017 NotPetya and Bad Rabbit cyberattacks.\n\n[All Posts by Cybereason Nocturnus](https://www.cybereason.com/blog/authors/cybereason-nocturnus)\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2021/2021-01-26 - Cybereason vs. RansomEXX Ransomware.pdf"
    ],
    "report_names": [
        "2021-01-26 - Cybereason vs. RansomEXX Ransomware.pdf"
    ],
    "threat_actors": [],
    "ts_created_at": 1673535762,
    "ts_updated_at": 1743041173,
    "ts_creation_date": 1653698824,
    "ts_modification_date": 1653698824,
    "files": {
        "pdf": "https://archive.orkl.eu/e9dc2ddf2cfda0ac448a34b02460fef4531c6c28.pdf",
        "text": "https://archive.orkl.eu/e9dc2ddf2cfda0ac448a34b02460fef4531c6c28.txt",
        "img": "https://archive.orkl.eu/e9dc2ddf2cfda0ac448a34b02460fef4531c6c28.jpg"
    }
}