# 0-raw-data/ORKL/library/jsons_entries/51ca09db-a6a4-430c-8986-dce4f149194e.json
{
    "id": "51ca09db-a6a4-430c-8986-dce4f149194e",
    "created_at": "2023-01-12T15:00:39.190752Z",
    "updated_at": "2025-03-27T02:05:44.266454Z",
    "deleted_at": null,
    "sha1_hash": "cde47d4a11beaf01ac642fba8bf2cc9f879ec1a9",
    "title": "2020-09-24 - Analysis Report (AR20-268A)- Federal Agency Compromised by Malicious Cyber Actor",
    "authors": "",
    "file_creation_date": "2022-05-29T01:05:28Z",
    "file_modification_date": "2022-05-29T01:05:28Z",
    "file_size": 332441,
    "plain_text": "# Federal Agency Compromised by Malicious Cyber Actor\n\n**us-cert.cisa.gov/ncas/analysis-reports/ar20-268a**\n\n## Summary\n\n_This Analysis Report uses the MITRE Adversarial Tactics, Techniques, and Common_\n_Knowledge (ATT&CK®) framework. See the_ _[ATT&CK for Enterprise framework for all](https://attack.mitre.org/versions/v7/matrices/enterprise/)_\n_referenced threat actor tactics and techniques._\n\nThe Cybersecurity and Infrastructure Security Agency (CISA) responded to a recent threat\nactor’s cyberattack on a federal agency’s enterprise network. By leveraging compromised\ncredentials, the cyber threat actor implanted sophisticated malware—including multi-stage\nmalware that evaded the affected agency’s anti-malware protection—and gained persistent\naccess through two reverse Socket Secure (SOCKS) proxies that exploited weaknesses in\nthe agency’s firewall.\n\n[For a downloadable copy of IOCs, see: AA20-268A.stix.](https://us-cert.cisa.gov/sites/default/files/publications/AR20-268A.stix.xml)\n\n## Description\n\nCISA became aware—via EINSTEIN, CISA’s intrusion detection system that monitors federal\ncivilian networks—of a potential compromise of a federal agency’s network. In coordination\nwith the affected agency, CISA conducted an incident response engagement, confirming\nmalicious activity. The following information is derived exclusively from the incident response\nengagement and provides the threat actor’s tactics, techniques, and procedures as well as\nindicators of compromise that CISA observed as part of the engagement.\n\n**Threat Actor Activity**\n\nThe cyber threat actor had valid access credentials for multiple users’ Microsoft Office 365\n(O365) accounts and domain administrator accounts, which they leveraged for Initial Access\n\n[[TA0001] to the agency's network (Valid Accounts](https://attack.mitre.org/versions/v7/tactics/TA0001/) [[T1078]). First the threat actor logged into](https://attack.mitre.org/versions/v7/techniques/T1078/)\na user’s O365 account from Internet Protocol (IP) address `91.219.236[.]166 and then`\nbrowsed pages on a SharePoint site and downloaded a file (Data from Information\n_Repositories: SharePoint_ [[T1213.002]). The cyber threat actor connected multiple times by](https://attack.mitre.org/versions/v7/techniques/T1213/002/)\nTransmission Control Protocol (TCP) from IP address `185.86.151[.]223 to the victim`\norganization’s virtual private network (VPN) server (Exploit Public-Facing Application\n\n[[T1190]).](https://attack.mitre.org/versions/v7/techniques/T1190/)\n\nCISA analysts were not able to determine how the cyber threat actor initially obtained the\ncredentials. It is possible the cyber actor obtained the credentials from an unpatched agency\nVPN server by exploiting a known vulnerability—CVE-2019-11510—in Pulse Secure\n(Exploitation for Credential Access [[T1212]). In April 2019, Pulse Secure released patches](https://attack.mitre.org/versions/v7/techniques/T1212/)\n\n\n-----\n\nfor several critical vulnerabilities—including CVE-2019-11510, which allows the remote,\n[unauthenticated retrieval of files, including passwords.[1] CISA has observed wide](https://kb.pulsesecure.net/articles/Pulse_Security_Advisories/SA44101)\n[exploitation of CVE-2019-11510 across the federal government.[2]](https://www.us-cert.gov/ncas/alerts/aa20-010a)\n\nAfter initial access, the threat actor performed Discovery [[TA0007] by logging into an agency](https://attack.mitre.org/versions/v7/tactics/TA0007/)\nO365 email account from `91.219.236[.]166 and viewing and downloading help desk`\nemail attachments with “Intranet access” and “VPN passwords” in the subject line, despite\nalready having privileged access (Email Collection [[T1114], Unsecured Credentials:](https://attack.mitre.org/versions/v7/techniques/T1114/)\n_Credentials In Files_ [[T1552.001]). (Note: these emails did not contain any passwords.) The](https://attack.mitre.org/versions/v7/techniques/T1552/001/)\nactor logged into the same email account via Remote Desktop Protocol (RDP) from IP\naddress `207.220.1[.]3 (External Remote Services` [[T1133]). The actor enumerated the](https://attack.mitre.org/versions/v7/techniques/T1133/)\nActive Directory and Group Policy key and changed a registry key for the Group Policy\n(Account Manipulation [[T1098]). Immediately afterward, the threat actor used common](https://attack.mitre.org/versions/v7/techniques/T1098/)\nMicrosoft Windows command line processes— conhost, `ipconfig,` `net,` `query,`\n```\nnetstat, ping, and whoami, plink.exe —to enumerate the compromised system and\n\n```\nnetwork (Command and Scripting Interpreter [[T1059], System Network Configuration](https://attack.mitre.org/versions/v7/techniques/T1059/)\n_Discovery_ [[T1016]).](https://attack.mitre.org/versions/v7/techniques/T1016/)\n\nThe cyber threat actor then attempted multiple times to connect to virtual private server\n(VPS) IP `185.86.151[.]223 through a Windows Server Message Block (SMB) client.`\nAlthough they connected and disconnected multiple times, the connections were ultimately\nsuccessful. During the same period, the actor used an alias secure identifier account they\nhad previously created to log into VPS `185.86.151[.]223 via an SMB share. The attacker`\nthen executed `plink.exe on a victim file server (Command and Scripting Interpreter`\n\n[[T1059]). ( plink.exe is a command-line version of PuTTy that is used for remote](https://attack.mitre.org/versions/v7/techniques/T1059/)\nadministration.)\n\nThe cyber threat actor established Persistence [[TA0003] and Command and Control](https://attack.mitre.org/versions/v7/tactics/TA0003/)\n\n[[TA0011] on the victim network by (1) creating a persistent Secure Socket Shell (SSH)](https://attack.mitre.org/versions/v7/tactics/TA0011/)\ntunnel/reverse SOCKS proxy, (2) running `inetinfo.exe (a unique, multi-stage malware`\nused to drop files), and (3) setting up a locally mounted remote share on IP address\n\n`78.27.70[.]237 (Proxy` [[T1090]). The mounted file share allowed the actor to freely move](https://attack.mitre.org/versions/v7/techniques/T1090/)\nduring its operations while leaving fewer artifacts for forensic analysis. Refer to Threat Actor\nMalware section for more information about the SSH Tunnel/reverse SOCKS proxy and\n```\ninetinfo.exe .\n\n```\nThe cyber threat actor created a local account, which they used for data Collection [[TA0009],](https://attack.mitre.org/versions/v7/tactics/TA0009/)\n_Exfiltration_ [[TA0010], Persistence](https://attack.mitre.org/versions/v7/tactics/TA0010/) [[TA0003], and Command and Control](https://attack.mitre.org/versions/v7/tactics/TA0003/) [[TA0011] (Create](https://attack.mitre.org/versions/v7/tactics/TA0011/)\n_Account_ [[T1136]). The cyber threat actor used the local account to:](https://attack.mitre.org/versions/v7/techniques/T1136/)\n\nBrowse directories on a victim file server (Data from Shared Network Drive [[T1039]).](https://attack.mitre.org/versions/v7/techniques/T1039/)\n\n\n-----\n\nCopy a file from a user s home directory to their locally mounted remote share (Data\n_Staged_ [[T1074]).](https://attack.mitre.org/versions/v7/techniques/T1074/)\n\nCISA analysts detected the cyber threat actor interacting with other files on users’\nhome directories but could not confirm whether they were exfiltrated.\nCreate a reverse SMB SOCKS proxy that allowed connection between an cyber threat\nactor-controlled VPS and the victim organization’s file server (refer to Threat Actor\nMalware section for more information) (Proxy [[T1090]).](https://attack.mitre.org/versions/v7/techniques/T1090/)\nInteract with PowerShell module `Invoke-TmpDavFS.psm (refer to Threat Actor`\nMalware section for more information).\nExfiltrate data from an account directory and file server directory using `tsclient`\n( tsclient is a Microsoft Windows Terminal Services client) (Data from Local System\n\n[[T1005], Data from Network Shared Drive](https://attack.mitre.org/versions/v7/techniques/T1005/) [[T1039]).](https://attack.mitre.org/versions/v7/techniques/T1039/)\nCreate two compressed Zip files with several files and directories on them (Archive\n_Collected Data_ [[T1560]); it is likely that the cyber threat actor exfiltrated these Zip files,](https://attack.mitre.org/versions/v7/techniques/T1560/)\nbut this cannot be confirmed because the actor masked their activity.\n\nSee figure 1 for the sequence of the cyber threat actor’s tactics and techniques.\n\n_Figure 1: Cyber threat actor tactics and techniques_\n\n**Threat Actor Malware**\n\nPersistent SSH Tunnel/Reverse SOCKS Proxy\n\nWhile logged in as “Administrator,” the cyber threat actor created two Scheduled Tasks (see\ntable 1) that worked in concert to establish a persistent SSH tunnel and reverse SOCKS\nproxy. The proxy allowed connections between an attacker-controlled remote server and one\n\n\n-----\n\nof the victim organization s file servers (Scheduled Task/Job [[T1053], Proxy](https://attack.mitre.org/versions/v7/techniques/T1053/) [[T1090]). The](https://attack.mitre.org/versions/v7/techniques/T1090/)\nReverse SOCKS Proxy communicated through port 8100 (Non-Standard Port [[T1571]). This](https://attack.mitre.org/versions/v7/techniques/T1571/)\nport is normally closed, but the attacker’s malware opened it.\n\n_Table 1: Scheduled Tasks composing SSH tunnel and reverse SOCKS proxy_\n\n**Scheduled Task** **Description**\n\n`ShellExperienceHost.exe` This task created a persistent SSH tunnel to attackercontrolled remote server `206.189.18[.]189 and`\nemployed port forwarding to allow connections from the\nremote server port 39999 to the victim file server through\nport 8100. This task was run daily.\n```\n                 ShellExperienceHost.exe is a version of plink.exe,\n\n```\na command-line version of PuTTy that is used for remote\nadministration.\n\n`WinDiag.exe` This task is a reverse SOCKS proxy that is preconfigured\nto bind to and listen on TCP port 8100. `WinDiag.exe`\nreceived responses through the SSH tunnel and\nforwarded the responses through port 8100 to the VPS IP\naddress `185.193.127[.]17 over port 443. This task`\nwas run on boot.\n```\n                 WinDiag.exe had compile information that matched the\n\n```\nVPS login name\n\n**Dropper Malware: inetinfo.exe**\n\nThe threat actor created a Scheduled Task to run inetinfo.exe (Scheduled Task/Job [[T1053]).](https://attack.mitre.org/versions/v7/techniques/T1053/)\n```\ninetinfo.exe is a unique, multi-stage malware used to drop files (figure 2). It dropped\nsystem.dll and 363691858 files and a second instance of inetinfo.exe . The\nsystem.dll from the second instance of inetinfo.exe decrypted 363691858 as\n\n```\nbinary from the first instance of `inetinfo.exe . The decrypted` `363691858 binary was`\ninjected into the second instance of `inetinfo.exe to create and connect to a locally`\nnamed tunnel. The injected binary then executed shellcode in memory that connected to IP\naddress `185.142.236[.]198, which resulted in download and execution of a payload.`\n\n\n-----\n\n_Figure 2: Dropper malware_ `inetinfo.exe`\n\nThe cyber threat actor was able to overcome the agency’s anti-malware protection, and\n```\ninetinfo.exe escaped quarantine. CISA analysts determined that the cyber threat actor\n\n```\naccessed the anti-malware product’s software license key and installation guide and then\nvisited a directory used by the product for temporary file analysis. After accessing this\ndirectory, the cyber threat actor was able to run `inetinfo.exe (Impair Defenses: Disable`\n_or Modify Tools_ [[T1562.001]).](https://attack.mitre.org/versions/v7/techniques/T1562/001/)\n\n**Reverse SMB SOCKS Proxy**\n\nPowerShell script `HardwareEnumeration.ps1 created a reverse SMB SOCKS proxy that`\nallowed connection between attacker-controlled VPS IP `185.193.127[.]18 and the victim`\norganization’s file server over port 443 (Command and Scripting Interpreter: Power Shell\n\n[[T1059.001], Proxy](https://attack.mitre.org/versions/v7/techniques/T1059/001/) [[T1090]). PowerShell script](https://attack.mitre.org/versions/v7/techniques/T1090/) `HardwareEnumeration.ps1 was executed`\ndaily via a Scheduled Task (Scheduled Task/Job [[T1053]).](https://attack.mitre.org/versions/v7/techniques/T1053/)\n```\nHardwareEnumeration.ps1 is a copy of Invoke-SocksProxy.ps1, a free tool created\n\n```\n[and distributed by a security researcher on GitHub.[3]](https://github.com/p3nt4/Invoke-SocksProxy) `Invoke-SocksProxy.ps1 creates a`\nreverse proxy from the local machine to attacker infrastructure through SMB TCP port 445\n(Non-Standard Port [[T1571]). The script was likely altered with the cyber threat actor’s](https://attack.mitre.org/versions/v7/techniques/T1571/)\nconfiguration needs.\n\n**PowerShell Module:** `invoke-TmpDavFS.psm`\n\n\n-----\n\n```\ninvoke-TmpDavFS.psm is a PowerShell module that creates a Web Distributed Authoring\n\n```\nand Versioning (WebDAV) server that can be mounted as a file system and communicates\nover TCP port 443 and TCP port 80. `invoke-TmpDavFS.psm is distributed on GitHub.[4]`\n\n## Solution\n\n**Indicators of Compromise**\n\nCISA analysts identified several IP addresses involved in the multiple stages of the outlined\nattack.\n```\n   185.86.151[.]223 – Command and Control (C2)\n   91.219.236[.]166 – C2\n   207.220.1[.]3 – C2\n   78.27.70[.]237 – Data Exfiltration\n   185.193.127[.]18 – Persistence\n\n```\n**Monitor Network Traffic for Unusual Activity**\n\nCISA recommends organizations monitor network traffic for the following unusual activity.\n\nUnusual open ports (e.g., port 8100)\nLarge outbound files\nUnexpected and unapproved protocols, especially outbound to the internet (e.g., SSH,\nSMB, RDP)\n\nIf network defenders note any of the above activity, they should investigate.\n\n## Prevention\n\nCISA recommends organizations implement the following recommendations to protect\nagainst activity identified in this report.\n\n**Deploy an Enterprise Firewall**\n\nOrganizations should deploy an enterprise firewall to control what is allowed in and out of\ntheir network.\n\nIf the organization chooses not to deploy an enterprise firewall, they should work with their\ninternet service provider to ensure the firewall is configured properly.\n\n**Block Unused Ports**\n\nOrganizations should conduct a survey of the traffic in and out of their enterprise to\ndetermine the ports needed for organizational functions. They should then configure their\nfirewall to block unnecessary ports. Organization should develop a change control process to\n\n\n-----\n\nmake control changes to those rules. Of special note, unused SMB, SSH, and FTP ports\nshould be blocked.\n\n**Additional Recommendations**\n\nCISA recommends organizations implement the following best practices.\n\nImplement multi-factor authentication, especially for privileged accounts.\nUse separate administrative accounts on separate administration workstations.\nImplement the principle of least privilege on data access.\nSecure RDP and other remote access solutions using multifactor authentication and\n“jump boxes” for access.\nDeploy and maintain endpoint defense tools on all endpoints.\nKeep software up to date.\n\n## References\n\n Revisions\n\nSeptember 24, 2020: Initial Version\n\n[This product is provided subject to this Notification and this](https://us-cert.cisa.gov/privacy/notification) [Privacy & Use policy.](https://www.dhs.gov/privacy-policy)\n\n**Please share your thoughts.**\n\n[We recently updated our anonymous product survey; we'd welcome your feedback.](https://www.surveymonkey.com/r/CISA-cyber-survey?product=https://us-cert.cisa.gov/ncas/analysis-reports/ar20-268a)\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2020/2020-09-24 - Analysis Report (AR20-268A)- Federal Agency Compromised by Malicious Cyber Actor.pdf"
    ],
    "report_names": [
        "2020-09-24 - Analysis Report (AR20-268A)- Federal Agency Compromised by Malicious Cyber Actor.pdf"
    ],
    "threat_actors": [],
    "ts_created_at": 1673535639,
    "ts_updated_at": 1743041144,
    "ts_creation_date": 1653786328,
    "ts_modification_date": 1653786328,
    "files": {
        "pdf": "https://archive.orkl.eu/cde47d4a11beaf01ac642fba8bf2cc9f879ec1a9.pdf",
        "text": "https://archive.orkl.eu/cde47d4a11beaf01ac642fba8bf2cc9f879ec1a9.txt",
        "img": "https://archive.orkl.eu/cde47d4a11beaf01ac642fba8bf2cc9f879ec1a9.jpg"
    }
}