# 0-raw-data/ORKL/library/jsons_entries/6d1b4aac-9b79-4c6d-a0b5-15287c63d509.json
{
    "id": "6d1b4aac-9b79-4c6d-a0b5-15287c63d509",
    "created_at": "2023-01-12T15:07:23.550713Z",
    "updated_at": "2025-03-27T02:05:38.465462Z",
    "deleted_at": null,
    "sha1_hash": "1cda0229788b6f25a11166ce3c29da05492bab2f",
    "title": "2021-06-17 - Klingon RAT Holding on for Dear Life",
    "authors": "",
    "file_creation_date": "2022-05-28T18:24:07Z",
    "file_modification_date": "2022-05-28T18:24:07Z",
    "file_size": 7172375,
    "plain_text": "# Klingon RAT Holding on for Dear Life\n\n**[intezer.com/blog/malware-analysis/klingon-rat-holding-on-for-dear-life/](https://www.intezer.com/blog/malware-analysis/klingon-rat-holding-on-for-dear-life/)**\n\n### Get Free Account\n\nJoin Now\n\n### Top Blogs\n\n**[Automate Alert Triage and Response Tasks with Intezer EDR Connect](https://www.intezer.com/blog/incident-response/alert-triage-edr-integrations/)**\n\n\nJune 17, 2021\n\n\nIntegrate with SentinelOne, CrowdStrike, and More One of the biggest pain points of cyber\nsecurity... [Read more](https://www.intezer.com/blog/incident-response/alert-triage-edr-integrations/)\n\n**[Elephant Framework Delivered in Phishing Attacks Against Ukrainian Organizations](https://www.intezer.com/blog/research/elephant-malware-targeting-ukrainian-orgs/)**\n\nA recently developed malware framework called Elephant is being delivered in targeted spear\nphishing campaigns... [Read more](https://www.intezer.com/blog/research/elephant-malware-targeting-ukrainian-orgs/)\n\n**[New Conversation Hijacking Campaign Delivering IcedID](https://www.intezer.com/blog/research/conversation-hijacking-campaign-delivering-icedid/)**\n\nThis post describes the technical analysis of a new campaign detected by Intezer’s research\nteam,... [Read more](https://www.intezer.com/blog/research/conversation-hijacking-campaign-delivering-icedid/)\n\n\n-----\n\n[With more malware written in Golang than ever before, the threat from Go-based Remote](https://www.intezer.com/resource/year-of-the-gopher-2020-go-malware-round-up/)\nAccess Trojans (RATs) has never been higher. Not only has the number of Go malware\nincreased but also the sophistication of these threats. This is a technical analysis of an\nadvanced RAT written in Go that we are calling Klingon RAT. The RAT is well-featured and\nresilient due to its multiple methods of persistence and privilege escalation. It was\ndetermined that the RAT is being used by cybercriminals for financial gain. It is important to\nstay on top of this threat as it will degrade Antivirus security through killing targeted processes\nand hiding communications through encrypted channels.\n\n## Technical Analysis\n\nWhen searching our various hunting platforms for malware one particular sample caught our\neye. This Go sample, active since at least 2019, was flagged as malicious but mostly unique\ncode by our platform. It is not common to find RATs with very few code reuse. Threat actors\nreuse code all the time to expedite malware development. Since it is rare to see a RAT with\nsuch a large amount of code written from scratch, we dug deeper down the gopher hole. This\nRAT is full of tactics to combat Antiviruses, maintain persistence and escalate privileges. It\ncommunicates encrypted with its Command and Control (C2) server using TLS and can\nreceive commands allowing the attacker to fully control the infected machine.\n\n_Figure 1: Old analysis with unique code_\n\n## Initialization\n\nThe malware starts by creating an object whose purpose is to store information about the\nvictim machine, controller setup and paths to dropped utilities. It will then run a WMI command\n(wmic process get Caption,ParentProcessId,ProcessId) to get all running processes. The\nreturned value is parsed and stored in a slice. The malware will check this process list and\nmatch it against a list of targeted Antivirus processes. The taskkill command is used to kill\n[matching processes and child processes. The targeted processes are linked here. To start](https://gist.github.com/aslefhewqiwbepqwefbpqsciwueh/a72fc5fee38a479e48d79710cd64d7c0)\ngathering the information on the victim machine, it will get the OS version using the\n[ver command, then grab the username. A GET request is made to https://api.ipify.org to get](https://api.ipify.org/)\nthe public IP address. Finally in this function, it will fetch the machine ID from the registry key\n\n\n-----\n\nHKEY_LOCAL_MACHINE\\SOFTWARE\\Microsoft\\Cryptography\\ as shown in Figure 2. This\nID will later be sent in a beacon to the Command and Control (C2) server.\n\n_Figure 2: Function that fetches the key_\n\n## Dependency Deployment\n\nThe malware will decompress and drop three Gzip embedded files into the %temp% directory.\nThe dropped files are utilities for the threat actor to use once a C2 channel has been\nestablished. The files dropped are Foxmail, PAExec and LSASS, shown below.\n\n\n-----\n\n_Figure 3: Head of embedded Foxmail.exe file, Gzip compressed_\n\n_Figure 4: Dropped dependencies Next, the malware will check to see if it is installed at “C:_\n\\Users\\IEUser\\AppData\\Local\\Windows Update\\updater10.exe.” If not installed, the malware\nwill be relocated to the path.\n\n## Persistence\n\nPersistence can be set up in multiple ways, some of which require admin privileges. Privilege\nescalation will be covered in a later section.\n\n## Registry Run Key: Current User\n\nThe following registry entry is created:\n\n**Key: Computer\\HKEY_CURRENT_USER\\Software\\Microsoft\\Windows\\CurrentVersion\\Run**\n**Name: Windows Updater**\n**Value: “C:\\Users\\\\AppData\\Local\\Windows Update\\updater10.exe” -1 -0**\n\n\n-----\n\n_Figure 5: Registry Run Key_\n\n## Registry Run Key: Local Machine\n\nA similar entry as the above is created at:\nComputer\\HKEY_LOCAL_MACHINE\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\n\n## Image File Execution Options Injection\n\nImage File Execution Options are configured by the Windows registry with the intention of\nbeing used for debugging. This can be leveraged for persistence as any executable can be\nused as a “debugger.” The malware ensures the following keys exist:\nHKEY_LOCAL_MACHINE Software\\Microsoft\\Windows NT\\CurrentVersion\\Accessibility\nHKEY_CURRENT_USER Software\\Microsoft\\Windows NT\\CurrentVersion\\Image File\nExecution Options\\magnify.exe The Image File Execution Options key has the following\nentries set:\n\n**Name** **Data**\n\nConfiguration mangnifierpane\n\nDebugger “C:\\Users\\IEUser\\AppData\\Local\\Windows Update\\updater10.exe” -1 -0\n\nThis causes the binary for Microsoft Screen Magnifier (magnify.exe) accessibility tool to be\nbackdoored and execute the malware.\n\n## WMI Event Subscription\n\nIn this option the malware utilizes “WMIC” to create an event subscription for persistence.\nThree commands are executed to create events in the “\\rootsubscription” namespace that will\nstart the payload within 60 seconds of Windows booting up. The commands executed are:\nwmic /namespace:’\\\\root\\subscription’ PATH __EventFilter CREATE\nName=’GuacBypassFilter’, EventNameSpace=’root\\cimv2′, QueryLanguage=’WQL’,\nQuery=’SELECT * FROM __InstanceModificationEvent WITHIN 60 WHERE TargetInstance\nISA ‘Win32_PerfFormattedData_PerfOS_System” wmic /namespace:’\\\\root\\subscription’\nPATH CommandLineEventConsumer CREATE Name=’GuacBypassConsumer’,\nExecutablePath='”C:\\Users\\IEUser\\AppData\\Local\\Windows Update\\updater10.exe” -1 -0′,\nCommandLineTemplate='”C:UsersIEUserAppDataLocalWindows Updateupdater10.exe” -1 -0′\n\n\n-----\n\nwmic /namespace: \\\\root\\subscription PATH __FilterToConsumerBinding CREATE\nFilter=’__EventFilter.Name=’GuacBypassFilter”,\nConsumer=’CommandLineEventConsumer.Name=’GuacBypassConsomer”)\n\n## Winlogon Helper DLL\n\nThe malware can modify the “Winlogon” key in order to run itself during Windows logon. The\npath of the executable is appended to the “Userinit” entry.\n\n_Figure 6: Winlogon registry modified_\n\n## Scheduled Task\n\nThe malware can create a scheduled task called “OneDriveUpdate” to maintain persistence.\nThe task is configured from an XML file, “elevator.xml” dropped to APPDATA, to trigger upon\nlogon.\n\n\n-----\n\n_Figure 7: Task configuration file_\n\n_Figure 8: Action of triggering the task The file “elevator.xml” is then removed from the disk._\n\n## Privilege Escalation\n\nThere are multiple avenues that the malware can take for privilege escalation. It will first test\nto see if it already has admin privileges and if it is a Windows server. To check if the process\nhas admin privileges, it will attempt to open “\\\\\\\\.\\\\PHYSICALDRIVE0;” if unsuccessful, the\nmalware will attempt to open “\\\\\\\\.\\\\SCSI0.” If successful for either of these, it will return “True”\nfrom the function. If “False,” the program will check to see if it is a Windows server by running\nthe command “systeminfo,” and parsing for the string “Microsoft Windows Server,” as shown\nin Figure 9.\n\n\n-----\n\n_Figure 9: Check for Windows Server The malware has four options for privilege escalation,_\none of which is not implemented properly:\n\n## UAC Bypass: Computer Defaults\n\nThis exploit starts by opening the following registry key: HKEY_CURRENT_USER\n(0x80000001) Software\\Classes\\ms-settings\\shell\\open\\command The default entry is set to\nthe path of the malware, and an entry “DelegateExecute” has an empty string value added.\nNext, the program “computerdefaults.exe” is executed to complete the exploit.\n\n_Figure 10: Registry set for exploit The key is deleted after exploitation._\n\n## UAC Bypass: Fodhelper\n\nThis exploit is similar to the Computer Defaults UAC bypass but this time it leverages the\nprogram “Features on Demand Helper” (Fodhelper.exe), a binary with the “autoelevate”\nsetting set to true. The same registry entries are used.\n\n\n-----\n\n_Figure 11: UAC bypass with Fodhelper.exe_\n\n## UAC Bypass: Disk Cleanup\n\nThis UAC bypass works by leveraging the scheduled task named “SilentCleanup.” This task\nruns with the highest privileges but is configured to have the ability to be executed by\nunprivileged users.\n\n\n-----\n\n_Figure 12: Config for SilentCleanup The malware attempts to leverage the environment_\nvariable “%windir%” to execute itself with higher privileges. The scheduled task runs an action\n“%windir%\\system32\\cleanmgr.exe,” therefore the malware tries to set the “windir” variable to\nthe path of the malware.\n\n_Figure 13: Action of the scheduled task (SilentCleanup)_\n\n_Figure 14: “windir” variable set in the registry After setting the registry, the malware runs the_\nscheduled task.\n\n\n-----\n\n_Figure 15: Execution of the scheduled task The resulting process:_\n\n_Figure 16: The elevated process_\n\n## UAC Bypass: Event Viewer\n\nBased on the strings in this path, it appears that the malware intended to leverage the “Event\nViewer” UAC bypass. But this does not appear to be properly implemented in the program.\n\n\n-----\n\n_Figure 17: References to “eventvwr” in a function called by “MakeAdmin” parent_\n\n## Command and Control\n\nBefore Command and Control (C2) is established the malware initiates a controller struct:\ntype control.Controller struct{     bot model.Bot     socksSessions []control.SocksProxy\nshellSessions []control.Shell     connection net.Conn     keepAlive net.Conn }\n\n\n-----\n\n[First, a x509 keypair is decoded from Base64 and loaded by the function tls.x509KeyPair.](https://golang.org/pkg/crypto/tls/#X509KeyPair)\n\n_[Figure 18: Loading x509 key pair The decoded keypair is linked here and](https://gist.github.com/aslefhewqiwbepqwefbpqsciwueh/b283988a0bf4f2d51f9a3f4697055901)_ [here. Strings from](https://gist.github.com/aslefhewqiwbepqwefbpqsciwueh/bd0ec3fbbf3b1bb1496dd8a53c4fb7e6)\nthis certificate can be matched to strings in the Issuer DN of a similar certificate with subject\n[“UrbanCulture, Inc.” A further PEM certificate is decoded and appended to the cert pool. A](https://search.censys.io/certificates/226a9e9c05935884da5507c7647c5e404faf349bdbaf170b9a4c064f7cc697fe)\nTLS handshake is performed with the C2 server 185.188.183[.]144 on the port 1141 and then\ncreates a Goroutine called “Controller.WaitCommands.” The malware is able to:\n\nStart a SOCKS proxy (‘proxy’)\nStart a reverse shell (‘shell’)\nStart an RDP server (‘rdp’)\nStart a binary (‘binary’)\nUpdate binary (‘update’)\nRun PowerShell command (‘cmd’)\n\nThe malware will initiate further Goroutines to collect information from the system. If running\nas administrator, it will run the Lsass binary previously dropped into the temp folder.\n\n_Figure 19: Path of the Lsass binary to be executed The results are stored in a file called_\n“Andrew.dmp” inside the temp folder. This information is sent to the C2 server through a HTTP\nPOST request.\n\n\n-----\n\n_Figure 20: Location of dump file Another routine will take a fingerprint of the machine,_\nconcatenating the results into a string, and send this off in a HTTP POST request. It runs the\nfollowing commands in this order:\n\n1. systeminfo\n2. ipconfig\n3. net view /all\n4. net view /all domain\n5. net users /domain\n6. nltest /domain_trusts\n7. nltest /domain_trusts /all_trusts\n\nFinally, the malware will periodically get information about the local network and adapters.\n\n## Detect and Respond to Klingon RAT\n\nDetect if your Windows machine or server has been compromised by Klingon RAT or any\n[variant that reuses code using the Intezer Analyze Live Endpoint Scanner available via the](https://analyze.intezer.com/endpoint-analyses)\n[enterprise edition. Running the scanner will classify all binary code residing in your machine’s](https://www.intezer.com/pricing/)\nmemory.\n\n_Figure 21: Endpoint scan of an infected system_\n\n## Indicators of Compromise\n\n\n-----\n\nMD5 C2\n\n8d44ccac6b5512a416339984ad664d79 185.188.183[.]144\n\n14471a353788bb6cdb6071d0e0a83004 94.177.123[.]134\n\n327090cbddf94fc901662f0e863ba0cb 88.214.27[.]40\n\n39d550fd902ca4c1461961d01ad1aeb6 51.83.216[.]211\n\n## MITRE ATT&CK\n\nTactic ID Name\n\n[Execution](https://attack.mitre.org/tactics/TA0002/) [T1059.001](https://attack.mitre.org/techniques/T1059/001/) PowerShell\n\n[T1059.003](https://attack.mitre.org/techniques/T1059/003/) Windows Command Shell\n\n[T1047](https://attack.mitre.org/techniques/T1047/) Windows Management Instrumentation\n\n[Persistence](https://attack.mitre.org/tactics/TA0003) [T1547.001](https://attack.mitre.org/techniques/T1547/001/) Registry Run Keys / Startup Folder\n\n[T1547.004](https://attack.mitre.org/techniques/T1547/004/) Winlogon Helper DLL\n\n[T1546.003](https://attack.mitre.org/techniques/T1546/003/) Windows Management Instrumentation Event\nSubscription\n\n[T1546.012](https://attack.mitre.org/techniques/T1546/012/) Image File Execution Options Injection\n\n[T1053.005](https://attack.mitre.org/techniques/T1053/005/) Scheduled Task\n\n[Privilege Escalation](https://attack.mitre.org/tactics/TA0004/) [T1548.002](https://attack.mitre.org/techniques/T1548/002/) Bypass User Account Control\n\n[Defense Evasion](https://attack.mitre.org/tactics/TA0005/) [T1562.001](https://attack.mitre.org/techniques/T1562/001/) Disable or Modify Tools\n\n[T1070.004](https://attack.mitre.org/techniques/T1070/004/) File Deletion\n\n[Credential Access](https://attack.mitre.org/tactics/TA0006/) [T1003.001](https://attack.mitre.org/techniques/T1003/001/) LSASS Memory\n\n[Discovery](https://attack.mitre.org/tactics/TA0007) [T1082](https://attack.mitre.org/techniques/T1082/) System Information Discovery\n\n[T1016](https://attack.mitre.org/techniques/T1016/) System Network Configuration Discovery\n\n[T1018](https://attack.mitre.org/techniques/T1018/) Remote System Discovery\n\n\nCommand and\nControl\n\n\n[T1571](https://attack.mitre.org/techniques/T1571/) Non-Standard Port\n\n[T1071.001](https://attack.mitre.org/techniques/T1071/001/) Web Protocols\n\n\n-----\n\n**Ryan Robinson**\nRyan is a security researcher analyzing malware and scripts. Formerly, he was a researcher\non Anomali's Threat Research Team.\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2021/2021-06-17 - Klingon RAT Holding on for Dear Life.pdf"
    ],
    "report_names": [
        "2021-06-17 - Klingon RAT Holding on for Dear Life.pdf"
    ],
    "threat_actors": [],
    "ts_created_at": 1673536043,
    "ts_updated_at": 1743041138,
    "ts_creation_date": 1653762247,
    "ts_modification_date": 1653762247,
    "files": {
        "pdf": "https://archive.orkl.eu/1cda0229788b6f25a11166ce3c29da05492bab2f.pdf",
        "text": "https://archive.orkl.eu/1cda0229788b6f25a11166ce3c29da05492bab2f.txt",
        "img": "https://archive.orkl.eu/1cda0229788b6f25a11166ce3c29da05492bab2f.jpg"
    }
}