# 0-raw-data/ORKL/library/jsons_entries/72335d06-57b2-4360-b832-0e9aed653f97.json
{
    "id": "72335d06-57b2-4360-b832-0e9aed653f97",
    "created_at": "2023-01-12T15:02:43.340904Z",
    "updated_at": "2025-03-27T02:05:40.210513Z",
    "deleted_at": null,
    "sha1_hash": "9b3963f27d6d51e73d97c152b2ea9ad8a2bc8fdd",
    "title": "2020-09-22 - MTR Casebook- Blocking a $15 million Maze ransomware attack",
    "authors": "",
    "file_creation_date": "2022-05-28T02:51:12Z",
    "file_modification_date": "2022-05-28T02:51:12Z",
    "file_size": 579930,
    "plain_text": "# MTR Casebook: Blocking a $15 million Maze ransomware attack\n\n**[news.sophos.com/en-us/2020/09/22/mtr-casebook-blocking-a-15-million-maze-ransomware-attack/](https://news.sophos.com/en-us/2020/09/22/mtr-casebook-blocking-a-15-million-maze-ransomware-attack/)**\n\nGreg Iddon September 22, 2020\n\n_Customer profile: An organization with many hundreds of networked devices based in Asia_\n_Pacific._\n\nThe Sophos Managed Threat Response (MTR) team was called in to help an organization\ntargeted with Maze ransomware. The attackers issued a ransom demand for US$15 million –\nif they had succeeded this would have been one of the most expense ransomware payments\nto date.\n\n## Background: Ransomware partners in crime\n\n[Maze is one of the most notorious ransomware families, active since 2019 when it evolved](https://news.sophos.com/en-us/2020/05/12/maze-ransomware-1-year-counting/)\nfrom ChaCha ransomware. It was among the first to combine data encryption with\ninformation theft.\n\nThe operators behind Maze have recently started colluding with other ransomware groups,\n[including LockBit,](https://news.sophos.com/en-us/2020/04/24/lockbit-ransomware-borrows-tricks-to-keep-up-with-revil-and-maze/) [SunCrypt and](https://www.bleepingcomputer.com/news/security/suncrypt-ransomware-sheds-light-on-the-maze-ransomware-cartel/) [Ragnar Locker, providing them with access to their platform](https://news.sophos.com/en-us/2020/05/21/ragnar-locker-ransomware-deploys-virtual-machine-to-dodge-security/)\nfor posting stolen victim data.\n\nThis appears to have led to a reciprocal sharing of tactics, techniques and procedures\n[(TTPs): in the attack covered here the Maze group borrowed a Ragnar Locker technique that](https://news.sophos.com/en-us/2020/05/21/ragnar-locker-ransomware-deploys-virtual-machine-to-dodge-security/)\ninvolves using virtual machines.\n\n\n-----\n\n**For detailed technical analysis of this collaboration between attackers read** **Maze**\n**attackers adopt Ragnar Locker virtual machine technique.**\n\n## Days 1-3: The attack begins\n\nPrior to the attack becoming active, the operators compromised a computer on the target’s\nnetwork. This computer was then used as a ‘beach head’ in the network. On multiple\noccasions during the attack, the attackers connected from here to other computers over\nRemote Desktop Protocol (RDP).\n\nOn day three, the main part of the attack began. The attackers exploited a domain admin\naccount with a weak password to take control of an unprotected Domain Controller (DC).\nThey then spent several days moving across the network.\n\n[Using the legitimate network scanning tool Advanced IP Scanner to map the network, the](https://www.advanced-ip-scanner.com/)\nattackers created lists of IP addresses to which they would later deploy ransomware. These\nincluded a list of the IP addresses of machines belonging to the target’s IT administrators.\n\nThe attackers’ attention then turned to the exfiltration of data.\n\nThey identified a file server and accessed it remotely over RDP using the compromised\ndomain admin account. Using the legitimate archiving tools WinRar and 7zip, they started\ncompressing folders located on it.\n\nThese archives were then copied back to the primary DC using the legitimate Total\nCommander FTP client that the attackers had installed on the file server.\n\n[The attackers tried to install the cloud storage application Mega on the DC. This was blocked](https://mega.nz/)\nas the target had added Mega to their blocked list using the application control capability in\n[Sophos Intercept X endpoint protection. The attackers then switched to using the web-based](https://www.sophos.com/en-us/products/endpoint-antivirus.aspx)\nversion instead, uploading the compressed files.\n\n## Days 4-5: The calm before the storm\n\nFor two days, the attackers went quiet. It’s likely they were waiting for a day when the\ntarget’s IT security team wouldn’t be working, like the weekend.\n\n## Day 6: The first ransomware attack is launched\n\nThe first Maze ransomware attack was launched on a Sunday, using the already\ncompromised domain admin account and the lists of IP addresses that had been identified.\n\nThis first attack actually comprised three attacks as the operators deployed three copies of\nthe Maze ransomware via batch scripts to the targeted computers:\n\n\n-----\n\nC:\\ProgramData\\enc6.exe\nC:\\ProgramData\\enc.exe\nC:\\ProgramData\\network.dll\n\nThree scheduled tasks were created to execute the ransomware:\n\n**Name** **Command**\n\nWindows Update Security Patches C:\\ProgramData\\enc6.exe\n\nWindows Update Security Patches 5 C:\\ProgramData\\enc.exe\n\nWindows Update Security regsvr32.exe /i c:\\programdata\\network.dll\n\nOver 700 computers were targeted in the attack, which was detected and blocked by Sophos\nIntercept X.\n\nEither the attackers didn’t realize the attack had been blocked or they were hoping that the\ntheft of the data would be enough for the target to pay up – but whatever the reason, upon\nlaunching the first attack attempt they issued a ransom demand for US$ 15 million.\n\n## Day 7: The MTR team gets to work\n\nRealizing that they were under attack, the target’s security team engaged the advanced\nincident response skills of the Sophos MTR team. Since they were not yet a Sophos MTR\ncustomer, the Sophos Rapid Response team was first engaged. The team quickly identified\nthe compromised admin account, identified and removed several malicious files, and blocked\nattacker commands and C2 (command and control) communications.\n\n## Day 8: Investigation and neutralization continue\n\nOver the following hours the MTR team found further tools and techniques used by the\nattackers, as well as evidence relating to the exfiltration of data. More files and accounts\nwere blocked.\n\n## Day 9: The second attack\n\nThe attackers launched a second attack via a different compromised account. This attack\nwas similar to the first one: commands were executed on a DC, looping through the lists of\nIP addresses contained in txt files.\n\nHowever, this time they coped a file called license.exe to C:\\ProgramData:\n\n\n-----\n\nThis was followed by a scheduled task to execute it. In this attack attempt the task was\ncalled “Google Chrome Security Update”:\n\nThe attack was quickly identified and stopped. Intercept X detected the ransomware, and the\nMTR team disabled and deleted both the compromised account and the license.exe file. No\nfiles were encrypted.\n\n## Day 9: Third time lucky?\n\nJust a few hours after the second attempt, the attackers tried again.\n\nBy now they seemed to be growing desperate. This attack targeted a single machine, the\nmain file server that the exfiltrated data had been taken from, and used a completely different\ntechnique to the previous attacks.\n\nIn the third attempt, the attackers distributed the ransomware payload inside a virtual\nmachine (VM).\n\nFortunately the MTR investigators recognized this new approach immediately as they had\n[also responded to the Ragnar Locker ransomware attack where the technique was first seen.](https://news.sophos.com/en-us/2020/05/21/ragnar-locker-ransomware-deploys-virtual-machine-to-dodge-security/)\n\nThe Maze operators had enhanced the technique, but it was undoubtedly the same. The\nattack was detected and stopped and no files were encrypted.\n\n## Defeating adversaries in human-led attacks\n\nThis casebook highlights how agile and adaptable human-operated attacks can be, with the\nattackers able to quickly substitute and reconfigure tools and return to the ring for another\nround. It also demonstrates how, to minimize likelihood of detection, attackers take\nadvantage of multiple legitimate IT tools in their attacks.\n\nSophos endpoint products detect components of this attack as Troj/Ransom[GAV or Troj/Swrort-EG. Indicators of compromise can be found on the SophosLabs Github.](https://github.com/sophoslabs/IoCs/blob/master/Ransomware-Maze.csv)\n\n## What can defenders do?\n\nThe most important things an IT security team can do is to reduce the attack surface,\nimplement strong security software, including specialist anti-ransomware security, educate\nemployees, and consider setting up or engaging a human threat hunting service to spot the\nclues that software can’t.\n\n\n-----\n\nAny organization can be a ransomware target, and any spam or phishing email, exposed\nRDP port, vulnerable exploitable gateway device or stolen remote access credentials will be\nenough for such adversaries to gain a foothold.\n\n## MITRE ATT&CK Mapping\n\nThe [MITRE ATT&CK framework is a globally accessible knowledge base of known adversary](https://attack.mitre.org/)\ntactics, techniques and procedures (TTPs). It can help security teams as well as threat\nhunters and analysts to better understand, anticipate and mitigate attacker behavior.\n\n**Initial Access**\n\n[T1078.002 – Valid Accounts: Domain Accounts](https://attack.mitre.org/techniques/T1078/002)\n[T1133 – External Remote Services](https://attack.mitre.org/techniques/T1133)\n\n**Execution**\n\n[T1059.001 – Command & Scripting Interrupter: PowerShell](https://attack.mitre.org/techniques/T1059/001)\n[T1059.003 – Command and Scripting Interpreter: Windows Command Shell](https://attack.mitre.org/techniques/T1059/003)\n[T1047 – Windows Management Instrumentation](https://attack.mitre.org/techniques/T1047)\n[T1053.005 – Scheduled Task/Job: Scheduled Task](https://attack.mitre.org/techniques/T1053/005)\n\n**Defense Evasion**\n\n[T1564.006 – Hide Artifacts: Run Virtual Instance](https://attack.mitre.org/techniques/T1564/006)\n\n**Credential Access**\n\n[T1003.003 – OS Credential Dumping](https://attack.mitre.org/techniques/T1003/003)\n\n**Discovery**\n\n[T1016 – System Network Configuration Discovery](https://attack.mitre.org/techniques/T1016)\n\n**Lateral Movement**\n\n[T1021.001 – Remote Services: Remote Desktop Protocol](https://attack.mitre.org/techniques/T1021/001)\n[T1021.002 – Remote Services: SMB/Windows Admin Shares](https://attack.mitre.org/techniques/T1021/002)\n\n**Command & Control**\n\n[T1071.001 – Application Layer Protocol: Web Protocols](https://attack.mitre.org/techniques/T1071/001)\n\n**Exfiltration**\n\n[T1567.002 – Exfiltration Over Web Service: Exfiltration to Cloud Storage](https://attack.mitre.org/techniques/T1567/002)\n\n\n-----\n\n**Impact**\n\n[T1486 – Data Encrypted for Impact](https://attack.mitre.org/techniques/T1486)\n\n## Sophos Managed Threat Response and threat hunting\n\n[For more information on the Sophos MTR service visit our website or](https://www.sophos.com/en-us/products/managed-threat-response.aspx) speak with a Sophos\nrepresentative.\n\n[If you prefer to conduct your own threat hunts Sophos EDR gives you the tools you need for](https://www.sophos.com/en-us/products/endpoint-antivirus/edr.aspx)\nadvanced threat hunting and IT security operations hygiene. Start a 30-day no obligation\ntrial today.\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2020/2020-09-22 - MTR Casebook- Blocking a $15 million Maze ransomware attack.pdf"
    ],
    "report_names": [
        "2020-09-22 - MTR Casebook- Blocking a $15 million Maze ransomware attack.pdf"
    ],
    "threat_actors": [],
    "ts_created_at": 1673535763,
    "ts_updated_at": 1743041140,
    "ts_creation_date": 1653706272,
    "ts_modification_date": 1653706272,
    "files": {
        "pdf": "https://archive.orkl.eu/9b3963f27d6d51e73d97c152b2ea9ad8a2bc8fdd.pdf",
        "text": "https://archive.orkl.eu/9b3963f27d6d51e73d97c152b2ea9ad8a2bc8fdd.txt",
        "img": "https://archive.orkl.eu/9b3963f27d6d51e73d97c152b2ea9ad8a2bc8fdd.jpg"
    }
}