# 0-raw-data/ORKL/library/jsons_entries/743fffc1-e9a1-45c5-b592-4b8648f21485.json
{
    "id": "743fffc1-e9a1-45c5-b592-4b8648f21485",
    "created_at": "2023-01-12T15:07:45.679974Z",
    "updated_at": "2025-03-27T02:06:04.991464Z",
    "deleted_at": null,
    "sha1_hash": "da779ccfeeadc22c7dda9b30ca0b1676bc3e4b2d",
    "title": "2021-03-05 - New steganography attack targets Azerbaijan",
    "authors": "",
    "file_creation_date": "2022-05-28T03:32:38Z",
    "file_modification_date": "2022-05-28T03:32:38Z",
    "file_size": 3945935,
    "plain_text": "# New steganography attack targets Azerbaijan\n\n**[blog.malwarebytes.com/threat-analysis/2021/03/new-steganography-attack-targets-azerbaijan/](https://blog.malwarebytes.com/threat-analysis/2021/03/new-steganography-attack-targets-azerbaijan/)**\n\nThreat Intelligence Team March 5, 2021\n\n_This blog post was authored by Hossein Jazi_\n\nThreat actors often vary their techniques to thwart security defenses and increase the efficiency of\ntheir attacks. One of the tricks they use is known as steganography and consists of hiding content\nwithin images.\n\nWe recently observed a malicious Word file that uses this technique to drop a Remote\nAdministration Trojan (RAT) that was new to us. Based on the decoy document, we assess that this\nattack is targeting the government and military of Azerbaijan.\n\nSince April 2020 attackers have been taking advantage of the tensions between Azerbaijan and\nArmenia to target Azerbaijanis. Researchers found several actors that have exploited this conflict\n[via phishing lures to drop AgentTesla and PoetRat. While AgentTesla has been distributed globally](https://blog.talosintelligence.com/2020/10/poetrat-update.html)\nthrough different spam campaigns, PoetRat has been used specifically to target Azerbaijanis.\n\nIt seems the document we analyze in this blog has no connection with PoetRat for several reasons:\nThe PoetRat actor has not used steganography in its malicious documents and used Python and\nLua variants while the actor we analyzed has dropped a .Net rat called Fairfax which does not\nseem to be a .Net variant of PoetRat.\n\n## Maldoc analysis\n\nThe document lure is written in Azerbaijani and talks about a “National Security and Scientific”\nconference that will be held in Azerbaijan in 2021.\n\n\n-----\n\nFigure 1: Maldoc\n\nlure content\nThe malicious document contains a macro that is obfuscated. The attacker has inserted random\ncharacters within the meaningful names to obfuscate the functions and variables names. Here are\nsome of the examples:\n\nAddArg_OACZT_20210214_115603_xokkn_uments29 -> AddArguments29\nzixokknpPath -> zipPath\n\n\n-----\n\ntesOACZTtcustomdirabcdefghijklmnopqrstuvwxyzect_OACZT_20210214_115603_xokkn_ory\n-> testcustomdirectory\n\nAfter deobfuscation, the names become clear and can easily figure out the intent of the macro.\n\nFigure 2: Macro after deobfuscation\nThe attacker also used another layer of obfuscation to encode strings. Function “MyFunc23” has\nbeen defined for this purpose. This function receives an array of numbers and decodes them into a\nstring.\n\nThis function has a loop that reads four numbers of the input array in each iteration and passes\nthem to another function to convert them to a character. At the end it concatenates those\ncharacters to build the final string.\n\nFigure\n\n3: decoder function\nThe convertor function defines a big switch statement that returned the character equivalent of\neach 4 numbers.\n\n\n-----\n\nFigure 4: Convertor function\nUpon opening the document and enabling the content the macro will be executed. At first it defines\nthe following files and directories:\n\nzipPath: Directory that stores the extracted zip file from png image\nappFolder: directory that stores the Rat\nrunner: path of the batch file which executes the Rat\ndocxPath: path of the file that keeps a copy of the current document\ndocxCopyPath: Path of the zip format of the copied document\ndocxUnzipFolder: Directory that contains the document after being unzipped\n\nFigure 5: Define names\n\n\n-----\n\nThen, it tries to create the appFolder directory and if it could not create the directory it exits. After\ncreating the directory, it copies itself in a new format to the file path defined before. The reason it\ncopies itself in a new format is because the current document is protected and even after unzipping\nits content the macro will not be able to find the image to extract the zip file.\n\nFigure 6: Creates a copy of itself\nTo create a copy of itself, It uses “SaveAs2” function that saves the specified document with a new\nname or format. The string “wdFormatDocumentDefault” has been passed in as a file format\nparameter which saves the document as DOCX format. In this way the macro can see the image\nthat has embedded zip file.\n\nFigure 7: Save as function\nIn the next step, it extracts the created document copy into the created folder and calls\n“ExtractFromPng” function to extract the embedded object from the png file. This function calls itself\nrecursively to read all chunk identifications within the png image until it reaches the “puNk” chunk\nidentification which is the chunk that has the embedded zip file. After finding the chunk, it extracts\nand writes it into “fairfax.zip”.\n\n\n-----\n\nFigure out 8: Extract zip file from png\n\n\n-----\n\nFigure 9: Chunk identification\nThe “fairfax.zip” is then extracted into %APPDATA%\\vstelmetry directory. It contains the an\nexecutable file (Fairfax.exe) as well as a batch file (runner.bat). The executable has been written in\nVisual Studio and it seems the attacker archived the whole Visual Studio project.\n\n\n-----\n\nFigure 10: Fairfax directory\nAt the end it performs some dummy functions and then executes runner.bat to execute fairfax.exe.\n\nFigure 11: Execute runner.bat\n\n## FairFax.exe:\n\n\n-----\n\nThis is a .Net RAT that has been developed using TAP model (Task Asynchronous Programming\nmodel). This model provides an abstraction over asynchronous code. In this model each\nfunctionality can be defined as a Task and will be executed based on the external resource\nallocation and when other tasks complete.\n\nFigure 12: Main\nThis RAT is not obfuscated and contains three main functionalities:\n\nDownload files\nUpload files\nTake screenshots\n\n\n-----\n\nFigure 13: Main functionalities\nAll the configurations have been stored in Global settings including appfoldername, vbfilename,\nhost address, scheduled task info, vbfile content, and cipherkey.\n\nFigure 14: Global settings\nAll the communications with the server are AES encrypted and base64 encoded.\n\n\n-----\n\nFigure 15: Decryption function\n\nFigure 16: Encryption Function\nFor network communications it has defined four different tasks to send and receive files and\ncommands: SendFileAsync, SendAsync, ReceiveAsync and ReceiveFileAsync.\n\n\n-----\n\nFigure 17:\n\nNetwork communications\nTo manage the files, it has FileManager class that can get the file and save into a temp directory\nand also zip files.\n\nFigure 18: File manager\nIt also has the capability to make itself persistent by creating a vbsfile and adding it to Scheduled\nTasks.\n\nFigure 19: Scheduled task\n\n\n-----\n\n## Conclusion\n\nThreat actors use many techniques to subvert analysis and detection; in this blog post we\nexamined a group employing the less common technique of steganography, in which the actor\nhides a malicious payload within an image.\n\nDue the geopolitical events happening between Azerbaijan and Armenia, digital attacks against\nthese countries have increased in the past year. Cisco Talos reported a new RAT named PoetRAT\nwhich was also used to target Azerbaijan, though differences in the sample analyzed in this post\nsuggest this RAT is not related. Malwarebytes analysts will continue to track this activity, and report\non any new findings related to this threat.\n\n## IOCs\n\n**File Name** **SHA256**\n\ntelebler.doc ef02527858797356c5e8571c5a22d00c481fbc9ce73c81a341d482ea3776878a\n\nFairfax.zip 4ad451a1c07d1760a0586c3c5132a68539d98198c402f4fc2b42b954ea9f76d7\n\n\nauroraXXXX.zip\n\nauroraXXXX.docx\n\n\n0573926b05c34af23c7003cc0a30cfc682335f7e787958f9be7e6804edacd0a1\n\n\n-----\n\nimage1.png f33db9011c69e6f4b13c765f77466de023f442d8a75bce8ab450f4978275671a\n\nrunner.bat 909a94451d2640f89ec25aebcede14f238ead06b94f28544a99f4ecc2411b3b5\n\nFairfax.exe ab0f4d290f3d4532896dea80563e342c825b12e0111c2d54eac62b1b942b854b\n\nFairfax.exe 69e880b0545330b8e6d1543c47d89b4907fb79899b40c2478c591225ffc551ce\n\n**C2:**\n\nvnedoprym.kozow[.]com\n\n111.90.150[.]37\n\n## MITRE ATT&CK Techniques\n\n**Tactic** **Id** **Name** **Details**\n\n**Initial Access** T1566 Phishing Distributing maldocs through\nphishing emails\n\n**Execution** T1059.003 Windows command shell Starts CMD.EXE for commands\nexecution\n\nT1064 Scripting Executing FairFax.exe using\nbatch file\n\nT1059.001 PowerShell Executes PowerShell scripts\n\nT1204.002 User Execution Manual execution by user\n\n**Persistence** T1053.005 Scheduled Task Uses Task Scheduler for\npersistence\n\n\n**Defense**\n**Evasion**\n\n\nT1140 Deobfuscate/Decode Files\nor Information\n\n\nThe RAT has the ability to decode\nbase\n\n64 data and decrypt AES\nencrypted data\n\n\n**Colletion** T1113 Screen Capture The RAT has the ability to capture\nthe screen\n\n\nThe RAT archived files using zip\nutility\n\nUsing HTTPS for C2\ncommunications\n\nC2 traffic are base64 encoded\nand AES encrypted\n\nExfiltrates the data over C2\n\n\n**Command and**\n**Control**\n\n\nT1560.001 Archive Collected\nData: Archive via Utility\n\nT1071.001 Application Layer\nProtocol: Web Protocols\n\nT1132.001 Data Encoding: Standard\nEncoding\n\n\n**Exfiltration** T1041 Exfiltration Over C2\nChannel\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2021/2021-03-05 - New steganography attack targets Azerbaijan.pdf"
    ],
    "report_names": [
        "2021-03-05 - New steganography attack targets Azerbaijan.pdf"
    ],
    "threat_actors": [],
    "ts_created_at": 1673536065,
    "ts_updated_at": 1743041164,
    "ts_creation_date": 1653708758,
    "ts_modification_date": 1653708758,
    "files": {
        "pdf": "https://archive.orkl.eu/da779ccfeeadc22c7dda9b30ca0b1676bc3e4b2d.pdf",
        "text": "https://archive.orkl.eu/da779ccfeeadc22c7dda9b30ca0b1676bc3e4b2d.txt",
        "img": "https://archive.orkl.eu/da779ccfeeadc22c7dda9b30ca0b1676bc3e4b2d.jpg"
    }
}