# 0-raw-data/ORKL/library/jsons_entries/8e42128c-3a00-4c8f-8db9-73be25bf5323.json
{
    "id": "8e42128c-3a00-4c8f-8db9-73be25bf5323",
    "created_at": "2023-06-05T02:06:24.819122Z",
    "updated_at": "2025-03-27T02:16:26.15213Z",
    "deleted_at": null,
    "sha1_hash": "dce176255da89448d986d6b13cbacf8f73b7c692",
    "title": "2023-04-18 - PrivateLoader- Analyzing the Encryption and Decryption of a Modern Loader",
    "authors": "",
    "file_creation_date": "2023-06-04T13:02:46Z",
    "file_modification_date": "2023-06-04T13:02:46Z",
    "file_size": 182797,
    "plain_text": "# PrivateLoader: Analyzing the Encryption and Decryption of a Modern Loader\n\n**[any.run/cybersecurity-blog/privateloader-analyzing-the-encryption-and-decryption-of-a-modern-loader/](https://any.run/cybersecurity-blog/privateloader-analyzing-the-encryption-and-decryption-of-a-modern-loader/)**\n\nkhr0x April 18, 2023\n\n[HomeMalware Analysis](https://any.run/cybersecurity-blog/)\nPrivateLoader: Analyzing the Encryption and Decryption of a Modern Loader\nIn this article, we delve into the inner workings of PrivateLoader, a notorious malicious loader family. We will explore the encryption and\ndecryption processes utilized by this malware, particularly focusing on its ability to protect itself using VMProtect, as well as its decryption of\nloaded libraries. Let’s dive in!\n\n\nPrivateLoader\n\n\n## PrivateLoader analysis introduction\n\nPrivateLoader is a malicious loader family, written in C++ and first discovered in early 2021.\n\nIt is known for distributing a wide range of malware, from simple information stealers to complex rootkits and spyware, utilizing payloads.\n\nThe distribution of this type of malware is managed by the Pay-Per-Install (PPI) service, a popular tool within the cybercriminal ecosystem that\ngenerates revenue by adding payloads to malware.\n\nThe code itself involves the decryption of loaded libraries.\nAt present, there are two versions of PrivateLoader available: one protected by VMProtect, and a regular version.\nEvery day, between 2 and 4 samples of this malware are uploaded.\n\n## Static Analysis of the Source File\n\n\n-----\n\n**SHA256: 27c1ed01c767f504642801a7e7a7de8d87dbc87dee88fbc5f6adb99f069afde4**\n\nUsing the Detect It Easy utility, we can see that the analyzed executable file is compiled in C++. There is no information about the packer,\nwhich could mean it was not possible to identify it.\n\n\nPrivateLoader's sample data\n\n\nFig. 1 – PrivateLoader’s sample data\n\n\nThe next step is to search for unencrypted strings using the strings command:\n```\nstrings --encoding=l loader.exe\n\n```\n\n-----\n\nInteresting strings detected in the executable file\n\n\nFig. 2 – Interesting strings detected in the executable file\nAnalyzing the discovered strings allows us to identify several interesting elements:\n\nA user-agent, which is likely used to masquerade as a legitimate browser application\nURL addresses for determining the current IP and geolocation\n\n## PrivateLoader dynamic analysis with ANY.RUN\n\n[We analyzed the sample in ANY.RUN interactive malware sandbox.](https://any.run/?utm_source=anyrunblog&utm_medium=article&utm_campaign=privateloader&utm_content=landing)\n\nHere’s a link to the task:\n\n[https://app.any.run/tasks/3e359dc7-934b-4ae1-89bf-ad33e346ed60](https://app.any.run/tasks/3e359dc7-934b-4ae1-89bf-ad33e346ed60?utm_source=anyrunblog&utm_medium=article&utm_campaign=privateloader&utm_content=task)\n\nThe process tree generated by the executable file appears as follows:\n\n\n-----\n\nPrivateLoader s process tree\n\n\nFig 3. – PrivateLoader’s process tree\n\n**Analyzing the process tree leads to the following conclusions:**\n\n1. The main PrivateLoader process creates a child process named “FhuC750omh76YtB1xgR7diEy.exe”, whose executable file is located in the\nuser’s “Pictures” directory (T1564 – Hide Artifacts):\n\nC:\\Users\\admin\\Pictures\\Minor Policy\n\n2. The created child process is added to the startup using Task Scheduler (T1053.005 – Scheduled Task/Job: Scheduled Task):\n\n**schtasks /create /f /RU “admin” /tr “”C:\\Program Files (x86)\\ClipManagerP0\\ClipManager_Svc.exe”” /tn “LOLPA4DESK HR” /sc HOURLY /rl**\nHIGHEST\n\nThe executable file of the child process was downloaded from the Internet (T1105 – Ingress Tool Transfer). We will not go into the detailed\nanalysis of it.\n\n\n-----\n\nPrivateLoader downloaded payload\n\n\nFig 4. – PrivateLoader downloaded payload\nAnalyzing the HTTP requests, we can observe connections and data exchanges with the C2 server (T1071.001 – Application Layer Protocol):\n\n\n-----\n\nPrivateLoader С2 addresses\n\n\nFig. 5 – С2 addresses\nThe content sent (as well as received) in POST requests consists of BASE64-encoded strings (T1132.001 – Data Encoding: Standard\nEncoding). Decoding these strings does not yield any readable results:\n\ndata=kSYhy9HPjD5Jhn9y6Evty4XFfJ3JgIwrSzln5bGnLfKDmbXix2ebDEXy6Ty3Bb8Hz2GB8w0Y2SL2JeBSZ4G80iHAkSS7JJyeiPwZOpWJONOFzEBar\nljR9hkvX_TJhqr1nNqQpYUB2lQ9i7NmmHeL_QSx8hUka_C3jOxi02ml5FyDDruXM_IWwPXvAGxtT8TVi9wLtfd0mF1O369GUAEeI45sF1pKeyDfssmqE=\n\nMoving forward to the indicators, we can see that the malware steals user credentials from browsers (T1552.001 Credentials In Files):\n\n\n-----\n\nPrivateLoader Stealing data\n\n\nFig. 6 – Stealing data\n\n## Technical Analysis of PrivateLoader \n\nFor the technical analysis, the following tasks were set:\n\n1. Locate the C2 server within the code\n\n1. Identify the encryption algorithms for the C2 server and, if possible, for strings as well.\n\n1. Automate the decryption of the C2 server and strings\n\nThe analysis of the executable file revealed that string encryption is done using the XOR algorithm (T1027 – Obfuscated Files or Information).\nInitially, the data and key are loaded into the stack, and then decrypted using the SIMD instruction “PXOR” and the “XMM” register. The result\nof the XOR operation is also stored in the stack.\n\n**The three stages of C2 server decryption are shown below.**\n\n1. Loading encrypted data into the stack:\n\n\n-----\n\nPrivateLoader data\n\n\nFig. 7 – Data\n\n1. Loading the encryption key into the stack:\n\n\n-----\n\nKey PrivateLoader\n\n\nFig. 8 – Key\n\n1. Decrypting the C2 server using the “PXOR” instruction and saving the results in the stack:\n\n\n-----\n\nDecrypting PrivateLoader\n\n\nFig. 9 – Decrypting\nDuring the analysis process, it was also found that the method similar to C2 decryption is used to decrypt the following:\n\nUsed API functions (T1027.007 – Obfuscated Files or Information: Dynamic API Resolution)\nPayloads\nURLs and more\n\nSome of the analyzed samples are protected by VMProtect. The search for string decryption is complicated by the fact that the decryption data\nis located in one function, while the XOR and key are in another. Moreover, the key is always the same.\n\n\n-----\n\nFig 10. – Decript VMprotect sample\n\n## Example of automating C2 server decryption of PrivateLoader\n\n[To automate the extraction of data and configuration, we can use the Triton framework. It will emulate code blocks that contain all the](https://github.com/JonathanSalwan/Triton)\nnecessary encrypted information.\n\nYou can find [an example of a script for emulating a specific block in our GitHub repository. The output of the script will be the decrypted C2](https://github.com/anyrun/blog-scripts/blob/main/Extractors/Privateloader/PrivateLoader.py?utm_source=anyrunblog&utm_medium=article&utm_campaign=privateloader&utm_content=script)\nserver.\n\n\n-----\n\nPrivateLoader Script output\n\n\nFig 11. – Script output\n\nTherefore, by emulating all the code blocks that contain encrypted data, we can obtain a set of strings with the necessary information,\nincluding the C2 server.\n\n## Extracting the PrivateLoader configuration\n\nIn our service, you can view the configuration, which is extracted automatically:\n\n\n-----\n\nPrivateLoader configution and strings\n\n\nFig. 12 –\n\nPrivateLoader configution and strings\nThe decrypted data includes C2 addresses and strings. The strings contain information such as: used libraries and their functions, registry\nkeys, paths to crypto wallets and browsers, etc.\n\n## Conclusion\n\nIn this article, we discussed encryption in PrivateLoader.\n\nIts main feature is the XOR of all strings it interacts with (C2, URLs, DLLs). Also, some samples are protected by VMprotect, which makes the\ncode a bit more complex due to the use of many functions.\n\n[If you’d like to read more content like this, read our LimeRAT Malware Analysis. Or check out our deep dive into the encryption and decryption](https://any.run/cybersecurity-blog/limerat-malware-analysis/)\nprocess of XLoader/FormBook.\n\n## MITRE (ARMATTACK)\n\n**Tactics** **Techniques** **Description**\n\n\nTA0007:\n\nSoftware discovery\n\n\nT1518:\n\nSoftware Discovery\n\nT1082:\n\nSystem Information\n\nDiscovery\n\n\nSearches for installed software\n\nin the system\n\nin the “Uninstall” key\n\nCollects system data\n\n\n-----\n\nTA0011:\n\nCommand and Control\n\n\nT1071.001:\n\nApplication Layer\n\nProtocol\n\n\nSending collected data\n\nto the control server\n\n\nT1105 Ingress Tool Transfer requests binary from the Internet\n\n\nT1132.001 – Data Encoding:\n\nStandard Encoding\n\n\nencode data with BASE64\n\n\nTA0006: Credential Access T1552.001: Credentials In Files Stealing of personal data – login data\n\nTA0005: Defense Evasion T1564 Hide Artifacts attempt to hide artifacts in user folder\n\n\nT1027.007 – Obfuscated Files or\n\nInformation: Dynamic\n\nAPI Resolution\n\n\nobfuscate then dynamically resolve API\n\nfunctions called by their malware\n\n\nT1027 – Obfuscated Files or Information attempt to make an executable or\n\nfile difficult to discover or\n\nanalyze by encrypting XOR\n\n\nTA0002: Execution T1053.005 – Scheduled\n\nTask/Job: Scheduled Task\n\n## IOCs\n\n**Title** **Description**\n\n\nabuse the Windows Task\n\nScheduler to create file in statup\n\n\nName 27c1ed01c767f504642801a7e7a7de8d87dbc87dee88fbc5f6adb99f069afde4 exe\n\nMD5 6cc7d9664c1a89c58549e57b5959bb38\n\nSHA1 85b665c501b9ab38710050e9a5c1b6d2e96acccc\n\nSHA256 27c1ed01c767f504642801a7e7a7de8d87dbc87dee88fbc5f6adb99f069afde4\n\n## Extracted URLs\n\nhttp://23[.]254[.]227[.]214/api/tracemap[.]php\nhttp://23[.]254[.]227[.]205/api/tracemap[.]php\nhttp://23[.]254[.]227[.]202/api/tracemap[.]php\n\nhttp://208[.]67[.]104[.]60/api/tracemap[.]php\nhttp://208[.]67[.]104[.]60/api/firegate[.]php\nhttp://163[.]123[.]143[.]4/download/YT_Client[.]exe\n\n## Dropped executable file\n\n**Title** **Description**\n\nName C:\\Users\\admin\\AppData\\Local\\Microsoft\\Windows\\INetCache\\IE\\AH8CR9J5\\YT_Client[1].exe\n\nSHA256 041f891934add72852c8fda245c95da959d7f98cc580383d198e42f2de039634\n\n## DNS requests\n\niplogger.org\nipinfo.io\nIplis.ru\n\n## Connections (IP)\n\n“23[.]254.227.214”\n“23[.]254.227.202”\n“23[.]254.227.205”\n\n“208[.]67.104.60”\n\n## MORE SAMPLES FOR YOUR RESEARCH\n\n\n-----\n\nttps //app a y u /tas s/ 8 [a6 6c](https://app.any.run/tasks/ff1872a6-6c1f-4f79-89da-995b9bd56152/) 9 89da 995b9bd56 5 /\n\n[https://app.any.run/tasks/6a8f93eb-be36-41bc-bf7f-534938a7e3a2/](https://app.any.run/tasks/6a8f93eb-be36-41bc-bf7f-534938a7e3a2/)\n\n[https://app.any.run/tasks/cc2cb367-82e9-4705-9767-8c12f7a67a21/](https://app.any.run/tasks/cc2cb367-82e9-4705-9767-8c12f7a67a21/)\n\n[https://app.any.run/tasks/c32312d8-4026-4a81-84e5-3d90ab2e309a/](https://app.any.run/tasks/c32312d8-4026-4a81-84e5-3d90ab2e309a/)\n\n[https://app.any.run/tasks/235754fa-6aa3-49dd-bbc4-1a7f9361f455/](https://app.any.run/tasks/235754fa-6aa3-49dd-bbc4-1a7f9361f455/)\n\nANY.RUN\nmalware analyst\n\n\nkhr0x\n\nI'm 21 years old and I work as a malware analyst for more than a year. I like finding out what kind of malware got on my computer. In my spare\ntime I do sports and play video games.\n\nUser avatar\nkhr0x\nMalware analyst at ANY.RUN\nI'm 21 years old and I work as a malware analyst for more than a year. I like finding out what kind of malware got on my computer. In my spare\ntime I do sports and play video games.\nWhat do you think about this post?\n\n12 answers\n\nAwful\nAverage\nGreat\n\nNo votes so far! Be the first to rate this post.\n\n0 comments\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2023/2023-04-18 - PrivateLoader- Analyzing the Encryption and Decryption of a Modern Loader.pdf"
    ],
    "report_names": [
        "2023-04-18 - PrivateLoader- Analyzing the Encryption and Decryption of a Modern Loader.pdf"
    ],
    "threat_actors": [
        {
            "id": "75108fc1-7f6a-450e-b024-10284f3f62bb",
            "created_at": "2024-11-01T02:00:52.756877Z",
            "updated_at": "2025-03-27T02:00:55.544216Z",
            "deleted_at": null,
            "main_name": "Play",
            "aliases": null,
            "source_name": "MITRE:Play",
            "tools": [
                "Nltest",
                "AdFind",
                "PsExec",
                "Wevtutil",
                "Cobalt Strike",
                "Playcrypt",
                "Mimikatz"
            ],
            "source_id": "MITRE",
            "reports": null
        }
    ],
    "ts_created_at": 1685930784,
    "ts_updated_at": 1743041786,
    "ts_creation_date": 1685883766,
    "ts_modification_date": 1685883766,
    "files": {
        "pdf": "https://archive.orkl.eu/dce176255da89448d986d6b13cbacf8f73b7c692.pdf",
        "text": "https://archive.orkl.eu/dce176255da89448d986d6b13cbacf8f73b7c692.txt",
        "img": "https://archive.orkl.eu/dce176255da89448d986d6b13cbacf8f73b7c692.jpg"
    }
}