# 0-raw-data/ORKL/library/jsons_entries/96369748-610b-4790-9e3b-d4b7f9251908.json
{
    "id": "96369748-610b-4790-9e3b-d4b7f9251908",
    "created_at": "2023-01-12T14:59:44.472347Z",
    "updated_at": "2025-03-27T02:09:18.48052Z",
    "deleted_at": null,
    "sha1_hash": "2a865e5045f4d6a7c752a9abdebb88148aec9287",
    "title": "2022-06-07 - Bumblebee Loader on The Rise",
    "authors": "",
    "file_creation_date": "2022-06-13T04:00:57Z",
    "file_modification_date": "2022-06-13T04:00:57Z",
    "file_size": 1651140,
    "plain_text": "# Bumblebee Loader on The Rise\n\n**blog.cyble.com/2022/06/07/bumblebee-loader-on-the-rise/**\n\n## Sophisticated loader delivers Cobalt-Strike Beacons\n\n\nJune 7, 2022\n\n\nIn March 2022, a new malware named “Bumblebee” was discovered and reportedly distributed via spam\ncampaigns. Researchers identified that Bumblebee is a replacement for BazarLoader malware, which has\ndelivered Conti Ransomware in the past. Bumblebee acts as a downloader and delivers known attack\nframeworks and open-source tools such as Cobalt Strike, Shellcode, Sliver, Meterpreter, etc. It also\ndownloads other types of malware such as ransomware, trojans, etc.\n\nOur intelligence indicates that the incidents of Bumblebee infection are on the rise, as shown below.\n\n\n-----\n\nRise of Bumblebee Malware\n\nThe Bumblebee infection starts through spam email. This email contains a link to further download an ISO\nfile that eventually drops the malicious Dynamic Link Library (DLL) file. The DLL file further loads\nBumblebee’s final payload on the victim’s machine.\n\nISO files are a type of archive file that contain an identical copy of data found on an optical disc, CDs,\nDVDs, etc. They are primarily used to back up optical discs or distribute large file sets intended to burn\nonto an optical disc.\n\nFigure 1 –\n\nBumblebee Infection Vector\n\n## Technical Details:\n\nThe complete technical analysis of Bumblebee is mentioned in the following sections. Cyble Research\nLabs analysed the hash (SHA256),\n“3e698d8d6e7820cc337d5e2eb3d8fbae752a4c05d11bcf00d3cb7d6dc45e1884” for analysis.\n\n### Bumblebee Initial Access:\n\nBumblebee has been distributed through spear-phishing email messages that use different methods to\ntrick users into downloading and opening the ISO files.\n\nThe spam email contains an HTML attachment as well as a hyperlink in the mail body to download the\nISO file. Similarly, the HTML attachment contains a link that downloads the ISO file from Microsoft\nOneDrive.\n\nFigure 2 shows the spam email that downloads ISO files from Microsoft OneDrive when users click on the\n“REVIEW THE DOCUMENT” hyperlink.\n\n\n-----\n\nFigure 2 –\n\nSpam Email (Source – Proofpoint)\nThe ISO file contains two files called Attachments.lnk and Attachments.dat. This malicious link file\ncontains the parameters to execute “Attachments.dat,” which is the Bumblebee payload, using Windows’\n_rundll32.exe service._\n\nFigure 3 shows the contents of the ISO file and properties of the .lnk file.\n\n\n-----\n\nFigure 3 – Contents of the ISO File\n\nand Properties of Malicious .lnk File\nTarget command line:\n\n_cmd.exe /c start /wait “” “C:\\Users\\Admin\\Local\\Temp\\Attachments.lnk” rundll32.exe_\n_“C:\\Windows\\System32\\rundll32.exe” Attachments.dat,ProcessLoad_\n\nIn another case of infection, the ISO file contains three files, namely New Folder.LNK, 7z.exe, and\n_arch.7z. The shortcut file New Folder.LNK launches powershell.exe and extracts the password-protected_\nfile arch.7z by using 7z.exe.\n\nThe arch.7z file contains a 64-bit DLL file named “arch.dll,” which is a Bumblebee loader. The PowerShell\ncommand extracts the arch.dll file into the location C:\\ProgramData\\ and executes it using rundll32.exe.\n\nFigure 4 shows the contents of the ISO file and properties of the .lnk file.\n\n\n-----\n\nFigure 4 – Contents of Malicious\n\nISO and Properties of .lnk file\n**Target command line:**\n\n_C:\\Windows\\System32\\cmd.exe /c powershell -WindowStyle Hidden -Command “.\\7za.exe x arch.7z_\n_-p434330cf2449 -o\\”c:\\programdata\\” -y > $null; rundll32 c:\\programdata\\arch.dll,oUlluzkNOs_\n\n### Defensive Evasion:\n\nBumblebee downloads and executes the other payloads on victim machines without being detected by\nany antivirus programs. Bumblebee uses various techniques to inject and attach the payloads into the\nrunning process.\n\nThe Bumblebee loader has a list of process names related to tools used by security researchers to identify\nif the malware is debugged or running in a virtual environment. The malware terminates its execution if it\nidentifies any of these processes running on the victim’s machine. The figure below shows the list of\nprocess names.\n\n\n-----\n\nFigure 5 – List of the Security Tools\n\nThe malware terminates its execution if it is identified to be running in a sandbox environment. The\nmalware calls the Wine_get_unix_file_name() API to identify the sandbox machine, as shown below.\n\nFigure 6 – Sandbox\n\nDetection using GetProcAddress()\nBumblebee also avoids running in the sandbox environment by comparing the victim’s specific usernames\nwith a list of hard-coded usernames. The hard-coded names are commonly-used sandbox usernames\nseen in the wild.\n\nIf user account names match with the names on the list, the malware terminates itself. A list of user\naccount names is shown in the figure below.\n\nFigure 7 – List of Hardcoded User Accounts\n\n\n-----\n\nThe malware performs additional checks to identify the virtual environment, such as Wine, Vbox, and\nVMware. To identify the virtual environment, the malware performs the following actions:\n\nQueries registry keys related to Virtual Machine-related software\nExecutes WMI queries to identify them\nIdentifies emulator by reading the respective registry keys\nIdentify the window name of the running process\n\nThis technique used by malware is highlighted in the figure below.\n\nFigure 8 –\n\nAdditional Defence Evasion Techniques\nAfter the evading detection, Bumblebee resolves its function names at runtime and creates a unique event\nname, 3C29FEA2-6FE8-4BF9-B98A-0E3442115F67.\n\nFigure 9 –\n\nBumblebee Creating Unique Event\nThe malware uses WMI queries to collect details such as system details, adapter details, etc., from the\nvictim’s machine. After that, it sends the stolen information to the Command and Control (C&C) server.\n\nFigure 10 – WMI Queries\n\n\n-----\n\nThe Bumblebee Loader uses various commands to perform malicious activities such as DLL injection,\ndownloading executables, uninstalling loaders, and enabling persistence. The commands used by the\nmalware are mentioned below.\n\n“dij”\n“dex”\n“sdl”\n“ins”\n\n### DLL Injection:\n\nThe malware receives the “dij” command for DLL and Shellcode injection. As shown in Figure 11, it injects\nShellcode into legitimate processes using the APC routine. It specifically injects code into the below\nprocesses:\n\n_\\\\Windows Photo Viewer\\\\ImagingDevices.exe_\n_\\\\Windows Mail\\\\wab.exe_\n_\\\\Windows Mail\\\\wabmig.exe_\n\nFigure 11\n\n– Process injection via Asynchronous Procedure Calls (APC)\nThe loader then creates two new sections within the target process and copies the Shellcode to the newly\ncreated sections to properly inject the Shellcode. Then it invokes the Shellcode in the target executable\nvia a dynamically resolved NtQueueApcThread().\n\n### Downloading Additional Payloads:\n\nThe malware receives the “dex” command for downloading and executing additional payloads. After\nreceiving this command along with payload data, it writes the file into a disk using the CreateFileA() and\n_WriteFile() functions and executes it via the COM object._\n\nIn this example, the malware uses the hardcoded name “wab.exe” to store the payload.\n\nFigure 12 – The dex command\n\noperation\n\n### Persistence:\n\n\n-----\n\nThe Ins command helps enable persistence by copying the Bumblebee malwareDLL into the %appdata%\ndirectory and creating a VBS script that loads the malicious DLL using a scheduled task.\nThe sdl command uses PowerShell to delete files from the infected system without prompting the user.\nThe PowerShell command used by the malware is:\n\n_PS C:\\> Remove-item -Path “filepath” -Force_\n\n### C&C Communication:\n\nThe figure below shows the COBALT STRIKE traffic from the malware.\n\nFigure 13 –\n\nCobalt Strike Network Traffic of Bumblebee Malware\n\n## Conclusion\n\nBumblebee is a new and highly sophisticated malware loader that employs extensive evasive maneuvers\nand anti-analysis tricks, including complex anti-virtualization techniques. To make the Bumblebee\nmalware’s activity stealthier and harder to detect, its Threat Actors frequently update these capabilities.\n\nBumblebee loader can be deployed to facilitate initial access and deliver payloads such as Cobalt Strike,\nransomware, etc. It is likely to become a popular tool for ransomware groups to deliver their payload.\n\nCyble Research Labs closely monitors the BumbleBee malware group and other similar Threat Actor\nactivities and analyzes them to better understand their motivations and keep our readers well-informed.\n\n## Our Recommendations\n\nRefrain from opening untrusted links and email attachments without first verifying their authenticity.\nEducate employees in terms of protecting themselves from threats like phishing’s/untrusted URLs.\nAvoid downloading files from unknown websites.\nUse strong passwords and enforce multi-factor authentication wherever possible.\nTurn on the automatic software update feature on your computer, mobile, and other connected\ndevices.\nUse a reputed antivirus and internet security software package on your connected devices, including\nPC, laptop, and mobile.\nBlock URLs that could spread the malware, e.g., Torrent/Warez.\nMonitor the beacon on the network level to block data exfiltration by malware or TAs.\nEnable Data Loss Prevention (DLP) Solutions on the employees’ systems.\n\n\n-----\n\n## MITRE ATT&CK® Techniques:\n\n**Tactic** **Technique ID** **Technique Name**\n\n\n**Initial Access** [T1566](https://attack.mitre.org/techniques/T1566/)\n[T1190](https://attack.mitre.org/techniques/T1190/)\n\n\nPhishing\nExploit Public-Facing Application.\n\n\n**Execution** [T1059](https://attack.mitre.org/techniques/T1059/) Command and Scripting Interpreter\n\n**Defence Evasion** [T1497](https://attack.mitre.org/techniques/T1497/) Virtualization/Sandbox Evasion\n\n**Persistence** [T1053](https://attack.mitre.org/techniques/T1053/) Scheduled Task/Job\n\n\n**Discovery** [T1012](https://attack.mitre.org/techniques/T1012/)\n[T1082](https://attack.mitre.org/techniques/T1082/)\n\n\nQuery Registry\nSystem Information Discovery\n\n\n**Credential Access** [T1552](https://attack.mitre.org/techniques/T1552/) Unsecured Credentials\n\n**Lateral Movement** [T1021](https://attack.mitre.org/techniques/T1021/) Remote Services\n\n**Impact** [T1496](https://attack.mitre.org/techniques/T1496/) Resource Hijacking\n\n## Indicators Of Compromise:\n\n**Indicators** **Indicator**\n\n**Type**\n\n\n**File name**\n\nNew-Folder00519.img\n\n7z.exe\n\narch.7z\n\nNew\nFolder.LNK\n\narch.dll\n\nLdrAddx64.dll\n\nwab.exe\n\n\n7092d2c4b041db8009962e865d6c5cd7\n\n11838141f869e74225be8bd0d4c866cb46ef0248\n\n0e859acbd03e59eae287b124803ec052cf027b519e608c7ccfd920044b9ee1c7\n\n42badc1d2f03a8b1e4875740d3d49336\ncee178da1fb05f99af7a3547093122893bd1eb46\nc136b1467d669a725478a6110ebaaab3cb88a3d389dfa688e06173c066b76fcf\n\n310803b7d4db43f2bd0040e21a4ef9fc\nf42c381524b5f52f0e1a5a8c60d62464b8644968\nb091415c1939d1da9a7d07901dd3d317a47b2a8ccc9c666d8cf53a512a80b8d6\n\nfd21be3db76b714cb4dfae779d1ada1f\n8157b198c00de0a19b1d02ae7b76c78857baccd2\n315b3d80643da454b40cc938a0e8794f90ccbd05868e55b4848cacbf047850ae\n\n16da4284ab7ab9d5669c34c339132ed6\n34dc625fc243d06cbc33d403ac7ee05edfd32819\n1249075a0c4af8ecfeb4a3ab1e9ef692cb8876591d73f3470106402ab1592717\n\nc9cf08565a10f4c46308037bd31a7f46\n17752edb2473b4a246d6a6980375bd87133e7514\n3e698d8d6e7820cc337d5e2eb3d8fbae752a4c05d11bcf00d3cb7d6dc45e1884\n\n33e03ca5dd9a8f85fdcf091a97312e45\n186981f889ad88a0d5f21c18adb8b35c78851c74\n64c299dc88a35d4ef551516be4f7ed95ae568a6ee0b66a1fcfc3f68bf80d87fe\n\n\nMD5\nSHA1\nSha256\n\nMD5\nSHA1\nSha256\n\nMD5\nSHA1\nSha256\n\nMD5\nSHA1\nSha256\n\nMD5\nSHA1\nSha256\n\nMD5\nSHA1\nSha256\n\nMD5\nSHA1\nSha256\n\n\n23.254.229[.]131 IP C&C\n\n\n-----\n\n79.110.52[.]71 IP C&C\n\n51.75.62[.]99 IP C&C\n\n23.106.215[.]123 IP C&C\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2022/2022-06-07 - Bumblebee Loader on The Rise.pdf"
    ],
    "report_names": [
        "2022-06-07 - Bumblebee Loader on The Rise.pdf"
    ],
    "threat_actors": [
        {
            "id": "610a7295-3139-4f34-8cec-b3da40add480",
            "created_at": "2023-01-06T13:46:38.608142Z",
            "updated_at": "2025-03-27T02:00:02.87217Z",
            "deleted_at": null,
            "main_name": "Cobalt",
            "aliases": [
                "Cobalt Gang",
                "GOLD KINGSWOOD",
                "COBALT SPIDER",
                "G0080",
                "Mule Libra",
                "Cobalt Group"
            ],
            "source_name": "MISPGALAXY:Cobalt",
            "tools": [],
            "source_id": "MISPGALAXY",
            "reports": null
        }
    ],
    "ts_created_at": 1673535584,
    "ts_updated_at": 1743041358,
    "ts_creation_date": 1655092857,
    "ts_modification_date": 1655092857,
    "files": {
        "pdf": "https://archive.orkl.eu/2a865e5045f4d6a7c752a9abdebb88148aec9287.pdf",
        "text": "https://archive.orkl.eu/2a865e5045f4d6a7c752a9abdebb88148aec9287.txt",
        "img": "https://archive.orkl.eu/2a865e5045f4d6a7c752a9abdebb88148aec9287.jpg"
    }
}