# 0-raw-data/ORKL/library/jsons_entries/a67c4381-0e13-4e92-b1bd-c9e1655e0a32.json
{
    "id": "a67c4381-0e13-4e92-b1bd-c9e1655e0a32",
    "created_at": "2023-02-02T02:07:56.19117Z",
    "updated_at": "2025-03-27T02:06:10.276639Z",
    "deleted_at": null,
    "sha1_hash": "93f0efc1ef9d1546601c9eabdcd5ae6c3fc4e3d0",
    "title": "2022-09-06 - Shikitega - New stealthy malware targeting Linux",
    "authors": "",
    "file_creation_date": "2023-02-01T07:50:16Z",
    "file_modification_date": "2023-02-01T07:50:16Z",
    "file_size": 1067660,
    "plain_text": "# Shikitega - New stealthy malware targeting Linux\n\n**[cybersecurity.att.com/blogs/labs-research/shikitega-new-stealthy-malware-targeting-linux](https://cybersecurity.att.com/blogs/labs-research/shikitega-new-stealthy-malware-targeting-linux)**\n\n1. [AT&T Cybersecurity](https://cybersecurity.att.com/)\n2. [Blog](https://cybersecurity.att.com/blogs)\n\nSeptember 6, 2022 | [Ofer Caspi](https://cybersecurity.att.com/blogs/author/ofer-caspi)\n\n## Executive summary\n\nAT&T Alien Labs has discovered a new malware targeting endpoints and IoT devices that are running Linux\noperating systems. Shikitega is delivered in a multistage infection chain where each module responds to a\npart of the payload and downloads and executes the next one. An attacker can gain full control of the\nsystem, in addition to the cryptocurrency miner that will be executed and set to persist.\n\n## Key takeaways:\n\nThe malware downloads and executes the Metasploit’s “Mettle” meterpreter to maximize its control on\ninfected machines.\nShikitega exploits system vulnerabilities to gain high privileges, persist and execute crypto miner.\nThe malware uses a polymorphic encoder to make it more difficult to detect by anti-virus engines.\nShikitega abuse legitimate cloud services to store some of its command and control servers (C&C).\n\n\n-----\n\nFigure 1. Shikitega operation process.\n\n## Background\n\n[With a rise of nearly 650% in malware and ransomware for Linux this year, reaching an all-time high in the](https://atlasvpn.com/blog/linux-malware-on-a-rise-reaching-all-time-high-in-h1-2022)\nfirst half year of 2022, threat actors find servers, endpoints and IoT devices based on Linux operating\nsystems more and more valuable and find new ways to deliver their malicious payloads. New malwares like\n[BotenaGo and](https://cybersecurity.att.com/blogs/labs-research/att-alien-labs-finds-new-golang-malwarebotenago-targeting-millions-of-routers-and-iot-devices-with-more-than-30-exploits) [EnemyBot are examples of how malware writers rapidly incorporate recently discovered](https://cybersecurity.att.com/blogs/labs-research/rapidly-evolving-iot-malware-enemybot-now-targeting-content-management-system-servers)\nvulnerabilities to find new victims and increase their reach.\n\nShikitega uses an infection chain in multiple layers, where the first one contains only a few hundred bytes,\nand each module is responsible for a specific task, from downloading and executing Metasploit meterpreter,\nexploiting Linux vulnerabilities, setting persistence in the infected machine to downloading and executing a\ncryptominer.\n\n## Analysis\n\nThe main dropper of the malware is a very small ELF file, where its total size is around only 370 bytes, while\nits actual code size is around 300 bytes. (figure 2)\n\n\n-----\n\nFigure 2. Malicious ELF file with a total of only 376 bytes.\n\n[The malware uses the “Shikata Ga Nai” polymorphic XOR additive feedback encoder, which is one of the](https://en.wikipedia.org/wiki/Shikata_ga_nai)\nmost popular encoders used in Metasploit. Using the encoder, the malware runs through several decode\nloops, where one loop decodes the next layer, until the final shellcode payload is decoded and executed.\nThe encoder stud is generated based on dynamic instruction substitution and dynamic block ordering. In\naddition, registers are selected dynamically. Below we can see how the encoder decrypts the first two\nloops: (figures 3 and 4)\n\nFigure 3. First “Shikata Ga Nai” decryption loop.\n\n\n-----\n\nFigure 4. Second “Shikata Ga Nai” decryption loop created by the first one.\n\nAfter several decryption loops, the final payload shellcode will be decrypted and executed. As the malware\ndoes not use any imports, it uses ‘int 0x80’ to execute the appropriate syscall. As the main dropper code is\nvery small, the malware will download and execute additional commands from its command and control by\ncalling 102 syscall (sys_socketcall). (Figure 5)\n\nFigure 5. Calling system functions using interrupts\n\nThe C&C will respond with additional shell commands to execute, as seen in the packet capture in figure 6.\nThe first bytes marked in blue are the shell commands that the malware will execute.\n\n\n-----\n\nFigure 6. Additional commands received from C&C.\n\nThe received command will download additional files from the server that won’t be stored in the hard drive,\nbut rather will be executed from memory only. (Figure 7)\n\nFigure 7. Executes additional shell code received from C&C.\n\n\n-----\n\nIn other malware versions, it will use the execve syscall to execute /bin/sh with command received from\nthe C&C. (figure 8)\n\nFigure 8. Executing shell commands by using syscall_execve.\n\n[The malware downloads and executes ‘Mettle’, a Metasploit meterpreter that allows the attacker to use a](https://github.com/rapid7/mettle)\nwide range of attacks from webcam control, sniffer, multiple reverse shells (tcp/http..), process control,\nexecute shell commands and more.\n\nIn addition the malware will use wget to download and execute the next stage dropper.\n\n## Next stage dropper\n\nThe next downloaded and executed file is an additional small ELF file (around 1kb) encoded with the\n“Shikata Ga Nai” encoder. The malware decrypts a shell command that will be executed by calling\nsyscall_execve with ‘/bin/sh” as a parameter with the decrypted shell. (Figure 9)\n\n\n-----\n\nFigure 9. Second stage dropper decrypts and executes shell commands.\n\nThe executed shell command will download and execute additional files. To execute the next and last stage\n[dropper, it will exploit two linux vulnerabilities to leverage privileges - CVE-2021-4034 and](https://cve.mitre.org/cgi-bin/cvename.cgi?name=CVE-2021-4034) [CVE-2021-3493](https://cve.mitre.org/cgi-bin/cvename.cgi?name=CVE-2021-3493)\n(figure 10 and 11).\n\nFigure 10. Exploiting Linux vulnerability CVE-2021-3493.\n\n\n-----\n\nFigure 11. Exploiting CVE-2021-4034 vulnerability.\n\nThe malware will leverage the exploit to download and execute the final stage with root privileges persistence and cryptominer payload.\n\n## Persistence\n\nTo achieve persistence, the malware will download and execute a total of 5 shell scripts. It persists in the\nsystem by setting 4 crontabs, two for the current logged in user and the other two for the user root. It will first\ncheck if the crontab command exists on the machine, and if not, the malware will install it and start the\ncrontab service.\n\n[To make sure only one instance is running, it will use the flock command with a lock file “/var/tmp/vm.lock”.](https://linux.die.net/man/2/flock)\n\nFigure 12. Adding root crontab to execute the final payload.\n\nBelow is the list of downloaded and executed script to achieve persistence:\n\n\nscript\nname\n\n\ndetails\n\n\nunix.sh Check if “crontab” commands exist in the system, if not install it and start the crontab\nservice.\n\nbrict.sh Adds crontab for current user to execute cryptominer.\n\npolitrict.sh Adds root crontab to execute cryptominer.\n\ntruct.sh Adds crontab for current user to download cryptominer and config from C&C.\n\n\n-----\n\nscript\nname\n\n\ndetails\n\n\nrestrict.sh Adds root crontab to download cryptominer and config from C&C.\n\nAs the malware persists with crontabs, it will delete all downloaded files from the system to hide its\npresence.\n\n## Cryptominer payload\n\nThe malware downloads and executes XMRig miner, a popular miner for the Monero cryptocurrency. It will\nalso set a crontab to download and execute the crypto miner and config from the C&C as mentioned in the\npersistence part above.\n\nFigure 13. XMRig miner is downloaded and executed on an infected machine.\n\n## Command and control\n\n[Shikitega uses cloud solutions to host some of its command and control servers (C&C) as shown by OTX in](https://otx.alienvault.com/indicator/hostname/dash.cloudflare.ovh)\nfigure 14. As the malware in some cases contacts the command and control server using directly the IP\nwithout domain name, it’s difficult to provide a complete list of indicators for detections since they are volatile\nand they will be used for legitimate purposes in a short period of time.\n\n\n-----\n\nFigure 14. Command and control server hosted on a legitimate cloud hosting service.\n\n## Recommended actions\n\n1. Keep software up to date with security updates.\n2. Install Antivirus and/or EDR in all endpoints.\n3. Use a backup system to backup server files.\n\n## Conclusion\n\nThreat actors continue to search for ways to deliver malware in new ways to stay under the radar and avoid\ndetection. Shiketega malware is delivered in a sophisticated way, it uses a polymorphic encoder, and it\ngradually delivers its payload where each step reveals only part of the total payload. In addition, the\nmalware abuses known hosting services to host its command and control servers. Stay safe!\n\n## Associated Indicators (IOCs)\n\nThe following technical indicators are associated with the reported intelligence. A list of indicators is also\navailable in the OTX Pulse. Please note, the pulse may include other activities related but out of the scope\nof the report.\n\nTYPE INDICATOR DESCRIPTION\n\nDOMAIN dash[.]cloudflare.ovh Command and\ncontrol\n\nDOMAIN main[.]cloudfronts.net Command and\ncontrol\n\n\n-----\n\nSHA256 b9db845097bbf1d2e3b2c0a4a7ca93b0dc80a8c9e8dbbc3d09ef77590c13d331 Malware hash\n\nSHA256 0233dcf6417ab33b48e7b54878893800d268b9b6e5ca6ad852693174226e3bed Malware hash\n\nSHA256 f7f105c0c669771daa6b469de9f99596647759d9dd16d0620be90005992128eb Malware hash\n\nSHA256 8462d0d14c4186978715ad5fa90cbb679c8ff7995bcefa6f9e11b16e5ad63732 Malware hash\n\nSHA256 d318e9f2086c3cf2a258e275f9c63929b4560744a504ced68622b2e0b3f56374 Malware hash\n\nSHA256 fc97a8992fa2fe3fd98afddcd03f2fc8f1502dd679a32d1348a9ed5b208c4765 Malware hash\n\nSHA256 e4a58509fea52a4917007b1cd1a87050b0109b50210c5d00e08ece1871af084d Malware hash\n\nSHA256 cbdd24ff70a363c1ec89708367e141ea2c141479cc4e3881dcd989eec859135d Malware hash\n\nSHA256 d5bd2b6b86ce14fbad5442a0211d4cb1d56b6c75f0b3d78ad8b8dd82483ff4f8 Malware hash\n\nSHA256 29aafbfd93c96b37866a89841752f29b55badba386840355b682b1853efafcb8 Malware hash\n\nSHA256 4ed78c4e90ca692f05189b80ce150f6337d237aaa846e0adf7d8097fcebacfe7 Malware hash\n\nSHA256 130888cb6930500cf65fc43522e2836d21529cab9291c8073873ad7a90c1fbc5 Malware hash\n\nSHA256 3ce8dfaedb3e87b2f0ad59e1c47b9b6791b99796d38edc3a72286f4b4e5dc098 Malware hash\n\nSHA256 6b514e9a30cbb4d6691dd0ebdeec73762a488884eb0f67f8594e07d356e3d275 Malware hash\n\nSHA256 7c70716a66db674e56f6e791fb73f6ce62ca1ddd8b8a51c74fc7a4ae6ad1b3ad Malware hash\n\nSHA256 2b305939d1069c7490b3539e2855ed7538c1a83eb2baca53e50e7ce1b3a165ab Malware hash\nCVE-20213493\n\nSHA256 4dcae1bddfc3e2cb98eae84e86fb58ec14ea6ef00778ac5974c4ec526d3da31f Malware hash\nCVE-20214034\n\nSHA256 e8e90f02705ecec9e73e3016b8b8fe915873ed0add87923bf4840831f807a4b4 Malware hash\n\nSHA256 64a31abd82af27487985a0c0f47946295b125e6d128819d1cbd0f6b62a95d6c4 Malware shell\nscript\n\n\n-----\n\nSHA256 623e7ad399c10f0025fba333a170887d0107bead29b60b07f5e93d26c9124955 Malware shell\nscript\n\nSHA256 59f0b03a9ccf8402e6392e07af29e2cfa1f08c0fc862825408dea6d00e3d91af Malware shell\nscript\n\nSHA256 9ca4fbfa2018fe334ca8f6519f1305c7fbe795af9eb62e9f58f09e858aab7338 Malware shell\nscript\n\nSHA256 05727581a43c61c5b71d959d0390d31985d7e3530c998194670a8d60e953e464 Malware shell\nscript\n\nSHA256 ea7d79f0ddb431684f63a901afc596af24898555200fc14cc2616e42ab95ea5d Malware hash\n\n## Mapped to MITRE ATT&CK\n\n[The findings of this report are mapped to the following MITRE ATT&CK Matrix techniques:](https://attack.mitre.org/)\n\nTA0002: Execution\n\nT1059: Command and Scripting Interpreter\nT1569: System Service\n\nT1569.002: Service Execution\nTA0003: Persistence\n\nT1543: Create or Modify System Process\nTA0005: Defense Evasion\n\nT1027: Obfuscated Files or Information\n\n### Share this with others\n\nTags: [malware research,](https://cybersecurity.att.com/blogs/tag/malware+research) [otx,](https://cybersecurity.att.com/blogs/tag/otx) [shikitega](https://cybersecurity.att.com/blogs/tag/shikitega)\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2022/2022-09-06 - Shikitega - New stealthy malware targeting Linux.pdf"
    ],
    "report_names": [
        "2022-09-06 - Shikitega - New stealthy malware targeting Linux.pdf"
    ],
    "threat_actors": [
        {
            "id": "aa73cd6a-868c-4ae4-a5b2-7cb2c5ad1e9d",
            "created_at": "2022-10-25T16:07:24.139848Z",
            "updated_at": "2025-03-27T02:02:10.120505Z",
            "deleted_at": null,
            "main_name": "Safe",
            "aliases": [],
            "source_name": "ETDA:Safe",
            "tools": [
                "DebugView",
                "LZ77",
                "OpenDoc",
                "SafeDisk",
                "TypeConfig",
                "UPXShell",
                "UsbDoc",
                "UsbExe"
            ],
            "source_id": "ETDA",
            "reports": null
        }
    ],
    "ts_created_at": 1675303676,
    "ts_updated_at": 1743041170,
    "ts_creation_date": 1675237816,
    "ts_modification_date": 1675237816,
    "files": {
        "pdf": "https://archive.orkl.eu/93f0efc1ef9d1546601c9eabdcd5ae6c3fc4e3d0.pdf",
        "text": "https://archive.orkl.eu/93f0efc1ef9d1546601c9eabdcd5ae6c3fc4e3d0.txt",
        "img": "https://archive.orkl.eu/93f0efc1ef9d1546601c9eabdcd5ae6c3fc4e3d0.jpg"
    }
}