# 0-raw-data/ORKL/library/jsons_entries/a99b669c-ff69-4076-b2b1-d3d9a10d7619.json
{
    "id": "a99b669c-ff69-4076-b2b1-d3d9a10d7619",
    "created_at": "2023-01-12T14:59:27.12606Z",
    "updated_at": "2025-03-27T02:05:30.152086Z",
    "deleted_at": null,
    "sha1_hash": "a8e3c1fb7a4e18e33f6b74203db0e5dbed8eaf27",
    "title": "2020-10-14 - Duck Hunting with Falcon Complete- Remediating a Fowl Banking Trojan, Part 3",
    "authors": "",
    "file_creation_date": "2022-05-28T22:07:21Z",
    "file_modification_date": "2022-05-28T22:07:21Z",
    "file_size": 1813198,
    "plain_text": "# Duck Hunting w/Falcon Complete Pt. 3: QakBot Countermeasures\n\n**[crowdstrike.com/blog/duck-hunting-with-falcon-complete-qakbot-countermeasures/](https://www.crowdstrike.com/blog/duck-hunting-with-falcon-complete-qakbot-countermeasures/)**\n\nThe Falcon Complete Team October 14, 2020\n\nThis blog is the last in a three-part series presenting the CrowdStrike® Falcon Complete™\nteam’s analysis of the recent QakBot campaigns observed in the wild and outlining a strategy for\n[the remote identification of a QakBot-infected host. While Part 1 and](https://www.crowdstrike.com/blog/duck-hunting-with-falcon-complete-analyzing-a-fowl-banking-trojan-part-1/) [Part 2 provided an analysis](https://www.crowdstrike.com/blog/duck-hunting-with-falcon-complete-qakbot-zip-based-campaign/)\nof techniques used by the threat actor to gain successful infections, Part 3 provides\nrecommendations for countermeasures that can be deployed via the CrowdStrike Falcon®\nplatform to prevent and contain infections before a widespread incident occurs. We also outline\na strategy for how the team helps organizations recover from incidents via our remote\nremediation capabilities.\n\nQakBot is an eCrime banking trojan that is capable of spreading laterally throughout a network.\nUtilizing a worm-like functionality, it spreads through brute forcing network shares, brute forcing\nActive Directory user group accounts or via SMB exploitation.\n\nQakBot also employs a robust set of anti-analysis features to evade detection and frustrate\n[analysis. Despite these protections, the CrowdStrike Falcon® platform detects and prevents the](https://www.crowdstrike.com/endpoint-security-products/falcon-platform/)\n[malware from completing its execution chain.](https://www.crowdstrike.com/epp-101/malware/)\n\n\n-----\n\n## Solution: Prevention, Containment and Remediation\n\nThe Falcon Complete team’s approach to defeating the QakBot threat can be defined in several\nways, depending on the tailored approach created for the customer. The best way to defend\nagainst QakBot is never allowing it to get a foothold in the first place, and the Falcon platform\noffers the prevention policies that are effective in stopping QakBot in its tracks. For situations\nwhere prevention is not enabled, a strategy of containment and fast remote response allows our\nanalysts to quickly remediate the infection via the Falcon Real Time Response (RTR) console,\nresulting in less downtime, less interference and more productivity.\n\n### Configuring Proper Prevention Policies \n\nPrevention policies are policies configured in the Falcon UI that allow organizations to customize\nhow aggressive the Falcon sensor is with detections and preventions. Organizations can choose\nto configure these policies themselves or take advantage of the Falcon Complete team’s\nexpertise to configure and offer guidance for best practices regarding these configurations.\n\nFigure 1. Prevention policy settings as shown in the Falcon UI (click image to enlarge)\n\nFalcon’s next-gen antivirus has the ability to block malware based on machine learning and\nbehavioral pattern analysis. Being a non-signature-reliant platform with the proper prevention\npolicies in place, Falcon reliably prevents known QakBot infections before any second-stage\npayloads are executed.\n\nWith prevention policies configured like those depicted in Figure 1, we can see that the\nexecution of a senate.m4a Zloader payload was blocked, and the detection details appear in the\nFalcon UI for analyst review (see Figure 2).\n\n\n-----\n\nFigure 2. Falcon blocking second-stage payload execution of senate.m4a (click image to enlarge)\n\n### Host Network Containment via the Falcon UI\n\nIn the event that prevention policies are not set to actively block the QakBot threat — as is\nsometimes the case in more sensitive environments — customers can take advantage of the\n[24/7/365 virtual security operations center (SOC) offered by Falcon Complete. In these cases,](https://www.crowdstrike.com/epp-101/security-operations-center-soc/)\nFalcon Complete will triage, stop lateral spread and remediate the QakBot banking trojan threat.\n\nIn the instance of an unprevented QakBot infection, the Falcon Complete team receives a highconfidence alert for malicious files, triggered by Falcon’s machine-learning algorithms.\n\nFigure 3. Detection in the UI as shown by Falcon (click image to enlarge)\n\nUpon expanding the alert, the analyst can recognize the QakBot threat by the tactics employed\n— many discovered during CrowdStrike’s tracking of this threat.\n\n\n-----\n\nFigure 4. Ability to network-contain hosts in emergency situations (click image to enlarge)\n\nIf among the pre-approved countermeasures, the analyst can network-contain the host to\nprevent the lateral spread of the info-stealer within the environment. This limits business impact\nand productivity loss for users, and saves time and cleanup for your internal SOC team.\n\n### Remediation with Falcon RTR\n\nQakBot malware is fairly simple in its functionality, but its capability to move laterally can be\npotentially devastating in enterprise networks. Because of this worm-like spreading and its costly\nrepercussions, the Falcon Complete team has classified the difficulty of the remediation process\nas “Hard.” The following is a brief illustration of the remote remediation process via Falcon Real\nTime Response (RTR).\n\n**The remediation of QakBot can be broken into three distinct steps:**\n\n1. Killing the malicious processes (e.g., injected explorer)\n\n\n-----\n\n2. Removing the persistence mechanism (e.g., scheduled task, registry run key)\n3. Removing disk artifacts (e.g., binaries and directories).\n\n**_Please note: Some of the examples in the following scenario have CrowdStrike Falcon_**\n**_configured with DETECTIONS ONLY and PREVENTIONS off for illustrative purposes. A_**\n**_properly configured Falcon instance, as noted previously, would prevent the activity_**\n**_presented here._**\n\n**STEP 1. Finding and Killing the Malicious** `explorer.exe Process`\n\nQakBot will create a new instance of the `explorer.exe process and` [inject itself into the new](https://attack.mitre.org/techniques/T1055/)\nprocess. It is possible for multiple instances of this process to occur, but responders must\ndetermine the malicious instance. It may not be immediately clear which `explorer.exe`\nprocess is legitimate, but this can be determined by querying the process with the “-Module”\nparameter, as shown in Figure 5.\n\nFigure 5. Query for explorer.exe process injection (click image to enlarge)\n\n[In the example in Figure 5, “gps” is a built-in alias for the “Get-Process” cmdlet, and the process](https://docs.microsoft.com/en-us/powershell/module/microsoft.powershell.management/get-process?view=powershell-7)\nid was determined by an early query (not shown). The legitimate explorer process will have\nseveral modules associated with it, but the malicious instance will have very few — typically less\nthan ten. In this case, Falcon had already prevented the injected process, but if there were a\nmalicious instance present, it could be killed with the following command: “kill ‘<PID>.’”\n\n**STEP 2. Removing Persistence**\n\n[QakBot typically employs a scheduled task and a registry run key for persistence mechanisms](https://attack.mitre.org/techniques/T1053/)\non compromised hosts. Figure 6 shows a query for scheduled tasks on the affected host. The\nsearch can be assisted by the historical knowledge that will point to a binary in %AppData%, so\nthis is used as a search string to quickly identify not only the malicious TaskName but also the\nfull path to the main binary.\n\n\n-----\n\nFigure 6. Query to find the path for scheduled task persistence mechanism (click image to enlarge)\n\nThe parent directory of the QakBot binary is a key indicator of compromise (IOC) and can be\nused later in the manual remediation process and as an integral search string. In the query\nshown in Figure 6, the parameter “-context 10” shows the surrounding ten lines and will reveal\nthe actual name of the task that is required for removal. This can simply be deleted with the\nbuilt-in CMD program schtasks.exe, like so:\n\nFigure 7. Removing the scheduled task (click image to enlarge)\n\nIn addition to scheduled tasks, QakBot will often create a registry run key to establish\npersistence. The modified value is placed under this key name:\nHKEY_CURRENT_USER\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run. Figure 8 shows\nthe output of the query for this registry key.\n\nFigure 8. Registry query (click image to enlarge)\n\nIn Figure 8, we obtain the value of the malicious run key and again encounter the path to the\nmain QakBot binary location. Removal of the key is performed with the command shown in\nFigure 9.\n\nFigure 9. Removing the run key (click image to enlarge)\n\n\n-----\n\n**STEP 3. Removing Remaining Artifacts**\n\nQakBot leaves file system residue in a few specific locations on infected hosts. The first one is\nthe random alphabetically named folder located in\n```\nC:\\Users\\*\\AppData\\Roaming\\Microsoft, identified in the remediation steps above. The\n\n```\nname of this directory is dynamic, but it will reliably contain QakBot’s core binary, .dat files and\nother resources. This initial download along with the temporary location of the loader’s initial\nexecution may also require removal, as shown in Figure 10.\n\nFigure 10. Removal of residual file system artifacts (click image to enlarge)\n\nThe steps outlined above are the general process for identifying QakBot artifacts and\nsuccessfully remediating a host for the QakBot malware. This process is quite effective in\nsingular instances, but in reality this is rarely the case due to QakBot’s lateral movement\ncapabilities. A properly configured Falcon platform is a critical component of a successful\ndefense strategy to defeat this threat.\n\n## Conclusion\n\nQakBot has undergone a resurgence in both its delivery volumes and technical evolution. The\npotential impact from a successful QakBot infection includes (but is certainly not limited to)\nwidespread lateral infections, theft of confidential data, deployment of secondary payloads and\nloss of organizational prestige.\n\nRemediating these types of infections becomes more complicated with these variants’ ability to\nspread laterally to many hosts. Furthermore, with a transition to a more remote workforce, the\ncapability to remotely remediate infections on hosts that are geographically distributed will\nbecome increasingly important as rebuilding systems becomes impractical.\n\nThe Falcon Complete team will continue to track this threat and monitor our clients’\nenvironments for any notable developments. Despite QakBot’s anti-analysis and evasive\ncapabilities, the CrowdStrike Falcon platform prevents this malware from completing its\nexecution chain when it detects the VBScript execution. The Falcon Complete team deals with\nthreats like this every day, providing our customers with the expertise required to remediate\nthese infections and help organizations recover from potentially devastating incidents.\n\n\n-----\n\n## Appendix\n\nTable 1 below contains a mapping of QakBot tactics to the MITRE ATT&CK framework.®\n\n\n**Tactic** **Technique** **Sub-Technique** **ID**\n\n\nInitial Access Phishing Spear-Phishing\nAttachment\n\nExecution User Execution Malicious Link, Malicious\nFile\n\n\n[T1566.001](https://attack.mitre.org/techniques/T1566/001/)\n\n[T1204.001,](https://attack.mitre.org/techniques/T1204/001) [T1204.002](https://attack.mitre.org/techniques/T1204/002)\n\n[T1059.001,](https://attack.mitre.org/techniques/T1059/001) [T1059.003,](https://attack.mitre.org/techniques/T1059/003)\n[T1059.005](https://attack.mitre.org/techniques/T1059/005)\n\n\nExecution Command and\nScripting Interpreter\n\nExecution Signed Binary Proxy\nExecution\n\nPersistence Boot or Logon\nAutostart Execution\n\n\nPowerShell, CMD Shell,\nVisual Basic\n\n\nRegistry Run Keys /\nStartup Folder\n\n\nMsiexec, Rundll32 [T1218.007,](https://attack.mitre.org/techniques/T1218/007/) [T1218.011](https://attack.mitre.org/techniques/T1218/011/)\n\n\n[T1547.001](https://attack.mitre.org/techniques/T1547/001)\n\n\nPersistence Scheduled Task/Job Scheduled Task [T1053.005](https://attack.mitre.org/techniques/T1053/005)\n\n\nDefense\nEvasion\n\nDefense\nEvasion\n\nDefense\nEvasion\n\n\nObfuscated Files or\nInformation\n\n\nVirtualization/Sandbox\nEvasion\n\n\nNone [T1027](https://attack.mitre.org/techniques/T1027/)\n\n\nProcess Injection Dynamic-link Library\nInjection\n\n\n[T1055.001](https://attack.mitre.org/techniques/T1055/001)\n\n\nSystem Checks [T1497.001](https://attack.mitre.org/techniques/T1497/001)\n\n\nDiscovery Virtualization/Sandbox\nEvasion\n\nDiscovery Network Share\nDiscovery\n\n\nUser Activity Based\nChecks\n\n\n[T1497.002](https://attack.mitre.org/techniques/T1497/002/)\n\n\nNone [T1135](https://attack.mitre.org/techniques/T1135/)\n\n\nCredential\nAccess\n\nLateral\nMovement\n\nCommand and\nControl\n\n\nBrute Force Password Guessing [T1110.001](https://attack.mitre.org/techniques/T1110/001/)\n\n\nApplication Layer\nProtocol\n\n\nRemote Services SMB/Windows Admin\nShares\n\n\n[T1021.002](https://attack.mitre.org/techniques/T1021/002/)\n\n\nWeb Protocols [T1071.001](https://attack.mitre.org/techniques/T1071/001/)\n\n\nTable 1. MITRE ATT&CK mapping\n\n**IOCs associated with QakBot analyses are available in Table 2.**\n\n**Indicator** **Purpose**\n\n\n-----\n\nPicturesViewer.dll, PicturesViewer.exe, PaintHelper.dll, PaintHelper.exe,\nfile1.exe\n\n\nQakBot binary\nnames\n\n\n“[0-9]{6,9}\\.zip”, “NUM_[0-9]{4,6}\\.vbs” Regular\nexpression of\nobserved\nfilename\nconvention of\nzip archives\ncontaining VBS\nthat launches\nQakBot\ndownloader\n\n“dfPEZd”, “ezQVN”, “wCdZgXH” Scheduled\nTask\ntasknames\n\n\n“””C:\\windows\\System32\\WScript.exe””\n“”C:\\Users\\*\\AppData\\Local\\Temp\\Temp1_*.zip\\NUM_*.vbs””\n\n\nCommand line\nexample of\ninitial\nexecution\n\n\nC:\\Users\\*\\Downloads\\“[0-9]{6,9}\\.zip” Initial QakBot\ndownload path.\nObserved as\nan 8 or 9character\nnumeric\nname.\n\nC:\\Users\\*\\AppData\\Local\\Temp\\Temp1_“[0-9]{6,9}\\.zip\\NUM_[0-9]{4,6}\\.vbs Execution path\nof VB\ndownloader\nscript\n\n\n%AppData%\\lwob\\esexydry.dll\n%AppData%\\PicturesViewer.dll\n\n%APPDATA%\\dasfdsfsdf.exe\n\n%APPDATA%\\Iwhoq\\pozypua.dll\n\n%APPDATA%\\IE\\GGYJG27Z\\dasfdsfs.df[1].exe\n\nC:\\Users\\Public\\tmpdir\n\n\nQakBot binary\npaths in home\ndirectories,\nobserved as a\nalphabetical\nname under an\nalphabetical\nfolder in\n%AppData%,\nor pre-named\nPicturesViewer,\nPaintHelper\n\n\nHKEY_CURRENT_USER\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run Registry run\nkey\npersistence\n\n\n-----\n\nTable 2. IOCs associated with QakBot\n\n**Additional Resources**\n\n_[Read “Duck Hunting with Falcon Complete: Analyzing a Fowl Banking Trojan, Part 1” and](https://www.crowdstrike.com/blog/duck-hunting-with-falcon-complete-analyzing-a-fowl-banking-trojan-part-1/)_\n_[“Duck Hunting with Falcon Complete: A Fowl Banking Trojan Evolves, Part 2” in this](https://www.crowdstrike.com/blog/duck-hunting-with-falcon-complete-qakbot-zip-based-campaign/)_\n_series._\n_Find out how CrowdStrike can help your organization answer its most important security_\n_questions:_ _[Visit the CrowdStrike Services webpage.](https://www.crowdstrike.com/services/)_\n_Learn how any size organization can achieve optimal security with Falcon Complete by_\n_visiting the product webpage._\n_[Learn more about Falcon X™ threat intelligence by visiting the webpage.](https://www.crowdstrike.com/endpoint-security-products/falcon-x-threat-intelligence/)_\n_Learn about CrowdStrike’s comprehensive next-gen endpoint protection platform by_\n_visiting_ _[the Falcon products webpage.](https://www.crowdstrike.com/endpoint-security-products/)_\n_Test CrowdStrike next-gen AV for yourself:_ _[Start your free trial of Falcon Prevent™.](https://go.crowdstrike.com/try-falcon-prevent.html)_\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2020/2020-10-14 - Duck Hunting with Falcon Complete- Remediating a Fowl Banking Trojan, Part 3.pdf"
    ],
    "report_names": [
        "2020-10-14 - Duck Hunting with Falcon Complete- Remediating a Fowl Banking Trojan, Part 3.pdf"
    ],
    "threat_actors": [],
    "ts_created_at": 1673535567,
    "ts_updated_at": 1743041130,
    "ts_creation_date": 1653775641,
    "ts_modification_date": 1653775641,
    "files": {
        "pdf": "https://archive.orkl.eu/a8e3c1fb7a4e18e33f6b74203db0e5dbed8eaf27.pdf",
        "text": "https://archive.orkl.eu/a8e3c1fb7a4e18e33f6b74203db0e5dbed8eaf27.txt",
        "img": "https://archive.orkl.eu/a8e3c1fb7a4e18e33f6b74203db0e5dbed8eaf27.jpg"
    }
}