# 0-raw-data/ORKL/library/jsons_entries/d05c8a8e-64db-4eb5-9fc3-ee4b2f48cbe7.json
{
    "id": "d05c8a8e-64db-4eb5-9fc3-ee4b2f48cbe7",
    "created_at": "2023-01-12T15:10:40.203029Z",
    "updated_at": "2025-03-27T02:05:42.081023Z",
    "deleted_at": null,
    "sha1_hash": "79b8086534717b2d50b80573ea3a588d4df4014a",
    "title": "2022-01-11 - TellYouThePass Ransomware Analysis Reveals a Modern Reinterpretation Using Golang",
    "authors": "",
    "file_creation_date": "2022-05-28T15:46:46Z",
    "file_modification_date": "2022-05-28T15:46:46Z",
    "file_size": 8517816,
    "plain_text": "# TellYouThePass Ransomware Analysis Reveals Modern Reinterpretation Using Golang\n\n**[crowdstrike.com/blog/tellyouthepass-ransomware-analysis-reveals-modern-reinterpretation-using-golang/](https://www.crowdstrike.com/blog/tellyouthepass-ransomware-analysis-reveals-modern-reinterpretation-using-golang/)**\n\nAnmol Maurya January 11, 2022\n\nTellYouThePass ransomware, discovered in 2019, recently re-emerged compiled using\nGolang\nGolang’s popularity among malware developers makes cross-platform development\nmore accessible\nTellYouThePass ransomware was recently associated with Log4Shell post-exploitation,\ntargeting Windows and Linux\nThe CrowdStrike Falcon® platform protects customers from Golang-written\nTellYouThePass ransomware using the power of machine learning and behavior-based\ndetection\n\n[The TellYouThePass ransomware family was recently reported as a post-exploitation](https://www.crowdstrike.com/cybersecurity-101/ransomware/)\nmalicious payload used in conjunction with a remote code execution vulnerability in Apache\nLog4j library, dubbed Log4Shell.\n\n\n-----\n\nTellYouThePass was first reported in early 2019 as a financially motivated ransomware\ndesigned to encrypt files and demand payment for restoring them. Targeting both Windows\nand Linux systems, TellYouThePass ransomware re-emerged in mid-December 2021 along\nwith other ransomware like Khonsari. This lesser-known ransomware family came back into\nthe spotlight as a post-exploitation payload associated with the Log4Shell. The remote code\nexecution vulnerability is estimated to expose affected organizations to a wave of\ncybersecurity risks.\n\nPreviously known TellYouThePass ransomware samples were written in traditional\nprogramming languages like Java or .Net., but two new recent samples reported in public\nrepositories have been rewritten and compiled in Golang.\n\n[Golang’s popularity among malware developers has steadily increased over the past years.](https://www.crowdstrike.com/blog/financial-motivation-drives-golang-malware-adoption/?&web_view=true)\nIt allows them to use the same codebase and compile it for all major operating systems,\nmaking cross-platform development work more accessible.\n\nWhat follows is a deeper dive into the new Golang-written TellYouThePass ransomware\nsamples for Windows and Linux and how the CrowdStrike Falcon platform protects against\nthem.\n\n## Setting Up the Analysis\n\nWe first check the binary for the “Go build id” string to identify the Golang build used for\ncompiling it. In recent campaigns of Go-written malware, especially in ransomware cases,\nattackers patch the binary to remove this string, making it difficult for researchers to use\nstring-based signatures to detect the binary as Go.\n\nGoing through the two samples —\n```\n460b096aaf535b0b8f0224da0f04c7f7997c62bf715839a8012c1e1154a38984 (Windows)\n5c8710638fad8eeac382b0323461892a3e1a8865da3625403769a4378622077e (Linux)\n\n```\n— we noticed that more than 85% of code in the Windows and Linux versions are almost the\nsame:\n\n\n-----\n\nFigure 1. The “main.” functions for both Windows and Linux samples are almost identical (Click to enlarge)\n\nA deeper dive into the some of the ransomware’s functions:\n\nFigure 2. TellYouThePass ransomware functions for the Windows sample in IDA Pro (Click to enlarge)\n\n[As we have previously discussed, we start by focusing on the “main.” functions in Golang.](https://www.crowdstrike.com/blog/financial-motivation-drives-golang-malware-adoption/)\nWe notice in this case that the malware authors have left only one main function and\nchanged the other functions to random names, making analysis difficult.\n\nThe sample checks the existence of the files “ showkey.txt ” and “ public.txt ” with the\nhelp of OS.Getenv, using “ ALLUSERSPROFILE “ and “ HOMEDRIVE “ as keys in Windows and\n```\nHome and /tmp/ in Linux. If it is present, it means encryption occurred, and it exists using\nruntime_gopanic ; otherwise, it creates them.\n\n```\n\n-----\n\nFigure 3. Encryption function followed by successful encryption for both Linux and Windows (Click to\n\nenlarge)\n\nFor Windows, the return is “ C:\\\\ProgramData “ and `/root/` directory in Linux. Using\n```\npath.join to join “ showkey.txt “ and “ public.txt ” with the directories results in:\n\n```\n**Windows** **Linux**\n\n\n” C:\\\\ProgramData/showkey.txt ”\n“ C:\\\\ProgramData/public.txt ”\n\n\n“ /root/showkey.txt ”\n“ /root/public.txt ”\n\n\nTable 1. Directories for saving showkey.txt and public.txt\n\n[The sample uses the Golang Crypto Packages for RSA key — some of them are](https://pkg.go.dev/crypto)\n[crypto_x509_MarshalPKCS1PublicKey,](https://pkg.go.dev/crypto/x509#MarshalPKCS1PublicKey) [crypto_x509_MarshalPKCS1PrivateKey,](https://pkg.go.dev/crypto/x509#MarshalPKCS1PrivateKey)\n[encoding_pem_EncodeToMemory and](https://pkg.go.dev/encoding/pem#EncodeToMemory) [crypto_rsa_GenerateMultiPrimeKey.](https://pkg.go.dev/crypto/rsa#GenerateMultiPrimeKey)\n\nAs seen in Figure 4, `crypto_x509_ MarshalPKCS1PrivateKey converts the RSA private`\nkey to PKCS #1, ASN.1 DER form. Then, [the encoding_pem_EncodeToMemory returns the](https://pkg.go.dev/encoding/pem#EncodeToMemory)\nPEM (Privacy Enhanced Mail) encoding, and after that, `runtime_slicebytetostring`\nconverts bytes to string, resulting in the conversion of bytes to string (see Figure 5).\n\n\n-----\n\nFigure 4. Function that generates the RSA private key\n\n\n-----\n\nFigure 5. The generated RSA key (Click to enlarge)\n\nThe RSA public key is generated using the\n```\nencoding_base64_ptr_Encoding_DecodeString and encoding_pem_encode packages\n\n```\nfrom Golang, as shown in Figure 6.\n\n\n-----\n\n-----\n\nFigure 6. Base64 decoding (Click to enlarge)\n\nAfter that, the PERSON_ID stores the encoding generated by\n[“encoding_base64__ptr_Encoding_EncodeToString” (in this case:](https://cs.opensource.google/go/go/+/go1.17.5:src/encoding/base64/base64.go;l=177)\n\n“ ABCDEFGHIJKLMNOPQRSTUVWXYZabcdefghijklmnopqrstuvwxyz0123456789 `+/ ”`\n\nas array for Base64 std encoding) every time the sample runs, saving it into\n“ showkey.txt ”.\n\nAfterward, another key is generated using the function below (Figure 7), also saving it into\n“ public.txt ”:\n\nFigure 7. Key generation function (Click to enlarge)\n\n## Ransomware Behavior Prior to Encryption\n\nTellYouThePass ransomware tries to kill some tasks and services before initiating the\nencryption routine, as shown in Table 2 below. However, in Linux, it requires root privilege to\ndo that. Targeted applications include various email clients, database applications, web\nservers and document editors.\n\nIt runs various commands using `cmd.exe` to kill tasks in Windows, and in Linux, it takes the\n[os_exec_command Go package to execute different commands using](https://cs.opensource.google/go/go/+/go1.17.5:src/os/exec/exec.go;l=169) `/bin/bash/ :`\n\n**Windows** **Linux**\n\n\n-----\n\n“taskkill /f /im msftesql.exe “\n“schtasks /delete /tn WM /F “\n“taskkill /f /im sqlagent.exe “\n“taskkill /f /im sqlbrowser.exe “\n“taskkill /f /im sqlservr.exe “\n“taskkill /f /im sqlwriter.exe “\n“taskkill /f /im oracle.exe “\n“taskkill /f /im ocssd.exe “\n“taskkill /f /im dbsnmp.exe “\n“taskkill /f /im synctime.exe “\n“taskkill /f /im mydesktopqos.exe “\n“taskkill /f /im agntsvc.exeisqlplussvc.”\n“taskkill /f /im xfssvccon.exe “\n“taskkill /f /im mydesktopservice.exe “\n“taskkill /f /im ocautoupds.exe “\n“taskkill /f /im agntsvc.exeagntsvc.exe “\n“taskkill /f /im agntsvc.exeencsvc.exe “\n“taskkill /f /im firefoxconfig.exe “\n“taskkill /f /im tbirdconfig.exe “\n“taskkill /f /im ocomm.exe “\n“taskkill /f /im mysqld.exe “\n“taskkill /f /im mysqld-nt.exe “\n“taskkill /f /im mysqld-opt.exe “\n“taskkill /f /im dbeng50.exe “\n“taskkill /f /im sqbcoreservice.exe “\n“taskkill /f /im excel.exe “\n“taskkill /f /im infopath.exe “\n“taskkill /f /im msaccess.exe “\n“taskkill /f /im mspub.exe “\n“taskkill /f /im onenote.exe “\n“taskkill /f /im outlook.exe “\n“taskkill /f /im powerpnt.exe “\n“taskkill /f /im steam.exe “\n“taskkill /f /im sqlservr.exe “\n“taskkill /f /im thebat.exe “\n“taskkill /f /im thebat64.exe “\n“taskkill /f /im thunderbird.exe “\n“taskkill /f /im visio.exe “\n“taskkill /f /im winword.exe “\n“taskkill /f /im wordpad.exe”\n“taskkill /f /im tnslsnr.exe”\n\n\n“service mysql stop”\n“/etc/init.d/mysqld stop”\n“service oracle stop”\n“systemctl disable \\”postgresql*\\””\n“systemctl disable \\”mysql*\\””\n“systemctl disable \\”oracle*\\””\n\n\nTable 2. TellYouThePass commands that try to terminate some tasks and services before\ninitiating the encryption routine\n\nAfter that, it iterates through all directories from A to Z and encrypts the files.\n\n\n-----\n\nBoth the Windows and the Linux versions have a list of directory exclusions for encryption,\nshown in Table 3.\n\n**Windows** **Linux**\n\n\nEFI.Boot\nEFI.Microsoft\nWindows\nProgram Files\nAll Users\nBoot\nIEidcache\nProgramData\ndesktop.ini\nautorun.inf\nnetuser.dat\niconcache.db\nthumbs.db\nLocal Settings\nbootfont.bin\nSystem Volume Information\nAppData\nRecycle.Bin\nRecovery\n\n\n/bin\n/boot\n/sbin\n/tmp\n/etc\n/lib\n/proc\n/dev\n/sys\n/usr/include\n/usr/java\n\n\nTable 3. TellYouThePass directory exclusions for encryption\n\nThe TellYouThePass ransomware focuses on encrypting popular media and file extensions,\nsaving their paths in the “ encfile.txt “ text file, located in the same folder as\n“ public.txt “ and “ showkey.txt ”.\n\nBelow is the full list of targeted extensions for encryption:\n\n1cd, 3dm, 3ds, 3fr, 3g2, 3gp, 3pr, 602, 7z, ps1, 7zip, aac, ab4, accdb, accde, accdr, accdt,\nach, acr, act, adb, adp, ads, aes, agdl, ai, aiff, ait, al, aoi, apj, arc, arw, asc, asf, asm, asp,\naspx, asx, avi, awg, back, backup, backupdb, bak, bank, bat, bay, bdb, bgt, bik, bin, bkp,\nblend, bmp, bpw, brd, c, cdf, cdr, cdr3, cdr4, cdr5, cdr6, cdrw, cdx, ce1, ce2, cer, cfg, cgm,\ncib, class, cls, cmd, cmt, conf, config, contact, cpi, cpp, cr2, craw, crt, crw, cs, csh, csl, csr,\ncss, csv, dac, dat, db, db3, db_journal, dbf, dbx, dc2, dch, dcr, dcs, ddd, ddoc, ddrw, dds,\nder, des, design, dgc, dif, dip, dit, djv, djvu, dng, doc, docb, docm, docx, dot, dotm, dotx, drf,\ndrw, dtd, dwg, dxb, dxf, dxg, edb, eml, eps, erbsql, erf, exf, fdb, ffd, fff, fh, fhd, fla, flac, flf, flv,\nflvv, fpx, frm, fxg, gif, gpg, gray, grey, groups, gry, gz, h, hbk, hdd, hpp, html, hwp, ibank, ibd,\nibz, idx, iif, iiq, incpas, indd, jar, java, jnt, jpe, jpeg, jpg, jsp, jspx, ashx, js, kc2, kdbx, kdc, key,\n\n\n-----\n\nkpdx, kwm, laccdb, lay, lay6, ldf, lit, log, lua, m, m2ts, m3u, m4p, m4u, m4v, mapimail, max,\nmbx, md, mdb, mdc, mdf, mef, mfw, mid, mkv, mlb, mml, mmw, mny, moneywell, mos, mov,\nmp3, mp4, mpeg, mpg, mrw, ms11, msg, myd, myi, nd, ndd, ndf, nef, nk2, nop, nrw, ns2,\nns3, ns4, nsd, nsf, nsg, nsh, nvram, nwb, nx2, nxl, nyf, oab, obj, odb, odc, odf, odg, odm,\nodp, ods, odt, ogg, oil, orf, ost, otg, oth, otp, ots, ott, p12, p7b, p7c, pab, pages, paq, pas,\npat, pcd, pct, pdb, pdd, pdf, pef, pem, pfx, php, pif, pl, plc, plus_muhd, png, pot, potm, potx,\nppam, pps, ppsm, ppsx, ppt, pptm, pptx, prf, ps, psafe3, psd, pspimage, pst, ptx, pwm, py,\nqba, qbb, qbm, qbr, qbw, qbx, qby, qcow, qcow2, qed, r3d, raf, rar, rat, raw, rb, rdb, rm, rtf,\nrvt, rw2, rwl, rwz, s3db, safe, sas7bdat, sav, save, say, sch, sd0, sda, sdf, sh, sldm, sldx, slk,\nsql, sqlite, sqlite3, sqlitedb, sr2, srf, srt, srw, st4, st5, st6, st7, so, st8, stc, std, sti, stm, stw,\nstx, svg, swf, sxc, sxd, sxg, sxi, sxm, sxw, tar, tar.bz2, tbk, tex, tga, tgz, thm, tif, tiff, tlg, txt,\nuop, uot, vb, vbox, vbs, vdi, vhd, vhdx, vmdk, vmsd, vmx, vmxf, vob, wab, wad, wallet, war,\nwav, wb2, wk1, wks, wma, wmv, wpd, wps, x11, x3f, xis, xla, xlam, xlc, xlk, xlm, xlr, xls, xlsb,\nxlsm, xlsx, xlt, xltm, xltx, xlw, xml, ycbcra, yuv, zip.\n\nFinally, the ransom note contains information about the encryption algorithm used to encrypt\nthe files, specifically RSA-1024 and AES-256. It also includes the personid, used for\nidentifying the victim. Following 0.05 bitcoin transfer into a designated and hardcoded wallet,\nattackers promise to provide victims with the decryption tool to recover all files.\n\nFigure 9. TellYouThePass ransom note (Click to enlarge)\n\n## CrowdStrike Falcon Protection\n\n\n-----\n\nThe Falcon platform automatically detects and protects against this type of Golang-written\nmalware using the power of the cloud, on-sensor and in-the-cloud machine learning, and\nindicators of attack (IOAs) to detect the threat. As Figure 10 shows, Falcon’s cloud-based\nmachine learning detects both Golang-written ransomware samples for TellYouThePass,\nimmediately protecting Windows and Linux environments.\n\nCrowdStrike Falcon leverages machine learning to identify known and unknown malware or\nthreats by understanding malicious intent. Both on-sensor and cloud-based machine learning\ncan detect and prevent post-exploitation threats leveraging exploits such as Log4Shell to\nprotect against malware, including the new Golang-written TellYouThePass ransomware.\n\nFigure 10. Falcon detection of Golang-written Windows TellYouThePass ransomware sample (Click to\n\nenlarge)\n\n\n-----\n\nFigure 11. Falcon detection of Golang-written Linux TellYouThePass ransomware sample (Click to\n\nenlarge)\n\nThe CrowdStrike Falcon platform provides protection against threats and visibility for all\nhosts in Windows, Linux and macOS, regardless of their location. The Falcon sensor can\ndetect and prevent threats ranging from ransomware, cryptocurrency miners, trojans and\nbotnets to stop today’s most sophisticated threats.\n\n## Indicators of Compromise (IOCs)\n\nFile/Host sha256\n\nWindows `460b096aaf535b0b8f0224da0f04c7f7997c62bf715839a8012c1e1154a38984`\n\nLinux `5c8710638fad8eeac382b0323461892a3e1a8865da3625403769a4378622077e`\n\n\nWindows\nhost\n\nLinux\nHost\n\n```\n45[.]76[.]99[.]222[:]80\n\n```\n```\n158[.]247[.]216[.]148[:]80\n\n```\n\n## MITRE ATT&CK® Framework Mapping\n\nAttack Id Tactic Description\n\nT1059 Execution Command and Scripting Interpreter\n\n\n-----\n\nT1053 Execution\nPersistence\n\nPrivilege Escalation\n\n\nScheduled Task/Job\n\n\nT1027 Defense Evasion Obfuscated Files or Information\n\nT1140 Defense Evasion Deobfuscate/Decode Files or Information\n\nT1083 Discovery File and Directory Discovery\n\nT1057 Discovery Process Discovery\n\nT1560 Collection Archive Collected Data\n\nT1486 Impact Data Encrypted for Impact\n\n**Additional Resources**\n\n_Read more about Golang malware in this blog: Golang Malware Is More than a Fad:_\n_Financial Motivation Drives Adoption_\n_Learn about another ransomware variant that uses a Golang packer: New_\n_Ransomware Variant Uses Golang Packer_\n_[Visit the product website to learn how the powerful CrowdStrike Falcon platform](https://www.crowdstrike.com/endpoint-security-products/falcon-platform/)_\n_provides comprehensive protection across your organization, workers and data,_\n_wherever they are located._\n_[Get a full-featured free trial of CrowdStrike Falcon Prevent™ and see how true next-](https://www.crowdstrike.com/resources/free-trials/try-falcon-prevent/)_\n_gen AV performs against today’s most sophisticated threats._\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2022/2022-01-11 - TellYouThePass Ransomware Analysis Reveals a Modern Reinterpretation Using Golang.pdf"
    ],
    "report_names": [
        "2022-01-11 - TellYouThePass Ransomware Analysis Reveals a Modern Reinterpretation Using Golang.pdf"
    ],
    "threat_actors": [
        {
            "id": "aa73cd6a-868c-4ae4-a5b2-7cb2c5ad1e9d",
            "created_at": "2022-10-25T16:07:24.139848Z",
            "updated_at": "2025-03-27T02:02:10.120505Z",
            "deleted_at": null,
            "main_name": "Safe",
            "aliases": [],
            "source_name": "ETDA:Safe",
            "tools": [
                "DebugView",
                "LZ77",
                "OpenDoc",
                "SafeDisk",
                "TypeConfig",
                "UPXShell",
                "UsbDoc",
                "UsbExe"
            ],
            "source_id": "ETDA",
            "reports": null
        }
    ],
    "ts_created_at": 1673536240,
    "ts_updated_at": 1743041142,
    "ts_creation_date": 1653752806,
    "ts_modification_date": 1653752806,
    "files": {
        "pdf": "https://archive.orkl.eu/79b8086534717b2d50b80573ea3a588d4df4014a.pdf",
        "text": "https://archive.orkl.eu/79b8086534717b2d50b80573ea3a588d4df4014a.txt",
        "img": "https://archive.orkl.eu/79b8086534717b2d50b80573ea3a588d4df4014a.jpg"
    }
}