# 0-raw-data/ORKL/library/jsons_entries/d7b39f24-529a-440a-82aa-9412e6c784b4.json
{
    "id": "d7b39f24-529a-440a-82aa-9412e6c784b4",
    "created_at": "2023-01-12T15:05:28.393577Z",
    "updated_at": "2025-03-27T02:05:27.380203Z",
    "deleted_at": null,
    "sha1_hash": "7be67ab8cfee0db3022a05e06953a1169e202c8f",
    "title": "2020-12-15 - QakBot reducing its on disk artifacts",
    "authors": "",
    "file_creation_date": "2022-05-28T04:52:32Z",
    "file_modification_date": "2022-05-28T04:52:32Z",
    "file_size": 957468,
    "plain_text": "# QakBot reducing its on disk artifacts\n\n**[hornetsecurity.com/en/threat-research/qakbot-reducing-its-on-disk-artifacts/](https://www.hornetsecurity.com/en/threat-research/qakbot-reducing-its-on-disk-artifacts/)**\n\nSecurity Lab December 15, 2020\n\n## Summary\n\nQakBot has been updated with more evasion techniques. QakBot’s configuration is now\nstored in a registry key instead of a file. The run key for persistence is not permanently\npresent in the registry but only written right before shutdown or reboot, and deleted\nimmediately after QakBot is executed again. QakBot’s executable is also not stored\npermanently on the file system anymore, but similarly to the run key registry entry, dropped\nonto the file system before reboots and deleted afterwards. This way security software can\nonly detect QakBot artifacts on disk, right before system shutdown, and shortly after system\nboot. However, at that time security software itself is shutting down and booting up, hence\nmay not detect QakBot’s new persistence method.\n\nOther changes include dynamic just-in-time decoding and destruction of strings at runtime.\nSo any string used in the malware is only decoded at runtime into memory only and\ndestroyed right afterwards.\n\nThe delivery method for the observed QakBot campaigns identified via the regular\nexpression pattern of `abc[0-9]+ is still XLM macro documents as reported previously.`\n\n## Background\n\nQakBot (also known as QBot, QuakBot, Pinkslipbot) has been around since 2008. It is\ndistributed via Emotet, i.e., Emotet will download QakBot onto victims that are already\ninfected with Emotet but it is also distributed directly via email. To this end, it uses email\nconversation thread hijacking in its campaigns,1 i.e., it will reply to emails that it finds in its\nvictim’s mailboxes. QakBot is known to escalate intrusions by downloading the ProLock\nransomware2 or lately the Egregor ransomware.\n\n\nThe observed QakBot campaigns identified by campaign ID `abc use XLM macro`\ndocuments for infection. We previously reported on their low detection.3\n\nAn overview of the current chain of infection used by the QakBot campaign with identifiers\nfollowing the regular expression pattern of `abc[0-9]+ can be seen in the following flow`\ngraph.\n\n\n-----\n\n## Technical Analysis\n\nIn the following analysis we briefly analyze the infection chain of QakBot after being\ndownloaded and launched by the malicious Excel document.\n\n\n-----\n\n### Evasion\n\nQakBot uses various evasion techniques to avoid detection by anti-virus software.\n\n**PE header manipulation**\n\nWe observed some QakBot DLLs with a manipulated PE header. The message text `This`\n```\nprogram cannot be run in DOS mode. has been altered.\n\n```\nThis seems like an attempt to circumvent some static detection rules matching for this\nmessage in the legacy MS-DOS stub of PE binaries.\n\n\n-----\n\n**Code signing**\n\nFirst, the initial downloaded and executed DLL is signed with a (at the time the analyzed\nsample was distributed) valid code signing certificate.\n```\n$ chktrust 904400.jpg\nMono CheckTrust - version 6.8.0.123\nVerify if an PE executable has a valid Authenticode(tm) signature\nCopyright 2002, 2003 Motus Technologies. Copyright 2004-2008 Novell. BSD licensed.\nWARNING! 904400.jpg is not timestamped!\nSUCCESS: 904400.jpg signature is valid\nand can be traced back to a trusted root!\n\n```\nThe signing CA is Sectigo and the organization is given as Aqua Direct s.r.o., which is an\nexisting company.\n\n\n-----\n\n```\n$ osslsigncode verify 904400.jpg \nCurrent PE checksum  : 00091021\nCalculated PE checksum: 00091021\nMessage digest algorithm : SHA1\nCurrent message digest  : 632DCB214EE9FB08441C640D240F672A7ABA6EB1\nCalculated message digest : 632DCB214EE9FB08441C640D240F672A7ABA6EB1\nSignature verification: ok\nNumber of signers: 1\n  Signer #0:\n    Subject: /C=CZ/postalCode=619 00/L=Brno/street=\\xC5\\xBDelezn\\xC3\\xA1\n646/8/O=Aqua Direct s.r.o./CN=Aqua Direct s.r.o.\n    Issuer : /C=GB/ST=Greater Manchester/L=Salford/O=Sectigo Limited/CN=Sectigo\nRSA Code Signing CA\nNumber of certificates: 4\n  Cert #0:\n    Subject: /C=CZ/postalCode=619 00/L=Brno/street=\\xC5\\xBDelezn\\xC3\\xA1\n646/8/O=Aqua Direct s.r.o./CN=Aqua Direct s.r.o.\n    Issuer : /C=GB/ST=Greater Manchester/L=Salford/O=Sectigo Limited/CN=Sectigo\nRSA Code Signing CA\n  Cert #1:\n    Subject: /C=GB/ST=Greater Manchester/L=Salford/O=Comodo CA Limited/CN=AAA\nCertificate Services\n    Issuer : /C=GB/ST=Greater Manchester/L=Salford/O=Comodo CA Limited/CN=AAA\nCertificate Services\n  Cert #2:\n    Subject: /C=US/ST=New Jersey/L=Jersey City/O=The USERTRUST\nNetwork/CN=USERTrust RSA Certification Authority\n    Issuer : /C=GB/ST=Greater Manchester/L=Salford/O=Comodo CA Limited/CN=AAA\nCertificate Services\n  Cert #3:\n    Subject: /C=GB/ST=Greater Manchester/L=Salford/O=Sectigo Limited/CN=Sectigo\nRSA Code Signing CA\n    Issuer : /C=US/ST=New Jersey/L=Jersey City/O=The USERTRUST\nNetwork/CN=USERTrust RSA Certification Authority\nSucceeded\n\n```\nIt is unknown whether the certificate was obtained from Sectigo by giving false information,\nthe certificate was stolen from Aqua Direct s.r.o., or whether the certificate was obtained from\nSectigo by giving stolen information from Aqua Direct s.r.o..\n\nQakBot is known to steal victim emails and use them in future malspam campaigns. So it is\nlikely that they also use stolen victim data to obtain code signing certificates. However, the\nactors behind QakBot can also buy the code signing certificate from a (malicious) third party.\n\n**Strings only decoded at runtime**\n\nQakBot will decode its strings only at runtime into memory. After usage the decoded strings\nare removed from memory again.\n\n\n-----\n\n**Processes**\n\nQakBot uses `CreateToolhelp32Snapshot and` `Process32{First,Next}W to enumerate`\nthe running processes.\n\nIt checks for the following processes:\n\n```\nWRSA exe\n\n```\n\n-----\n\n```\nvkise.exe\n\n```\n\nQakBot will set specific bits in a bit mask for each running process it finds. Depending on the\nresulting bit mask the further infection path is altered, e.g., if `avp.exe has been`\nencountered. QakBot will later inject its code into `mobsync.exe instead of` `explorer.exe .`\nBecause the searched process names are related to security solutions, we believe that this\nway QakBot tailors its execution path to evade detection by specific vendors.\n\nThen in another loop, again using `CreateToolhelp32Snapshot and`\n```\nProcess32{First,Next}W, it checks for:\n   srvpost.exe\n   frida-winjector-helper-32.exe\n   frida-winjector-helper-64.exe\n\n```\nIf it detects any of those processes the execution flow will run into a loop continuously calling\n```\nWaitForSingleObject(handle, 0x1fa) on a handle previously generated via\nCreateEvent(NULL, FALSE, FALSE, ...), i.e., it runs in an infinite loop.\n\n```\n**Device drivers**\n\nNext, QakBot uses `SetupDiGetDeviceRegistryPropertyA (querying properties`\n```\nSPDRP_DEVICEDESC and SPDRP_SERVICE ) to check for device drivers containing the\n\n```\nfollowing strings:\n```\n   VBoxVideo\n   Red Hat VirtIO\n   QEMU\n   A3E64E55_pr\n\n```\nWe believe the search for `A3E64E55_pr is used to detect an artifact of the ANY.RUN`\nsandbox.4 Alternatively, but unlikely, it could be used to detect an artifact of the long ago\ndefunct xCore Complex Protection AV solution using a similar driver with the name\n```\nA3E64E55_pr.sys .\n\n```\nIf it detects any of those device drivers the execution flow will run into the same infinite loop\ncontinuously calling `WaitForSingleObject(handle, 0x1fa) on a` `handle previously`\ngenerated via `CreateEvent(NULL, FALSE, FALSE, ...), as previously mentioned.`\n\n\n-----\n\n**Process injection**\n\nQakBot starts `C:\\Windows\\SysWOW64\\explorer.exe in suspended state and injects a DLL`\ninto it using `CreateProcessInternalW,` `NtMapViewOfSection,`\n```\nNtAllocateVirtualMemory, WriteProcessMemory, memcpy,\n\n```\n\n-----\n\n-----\n\n[The injected DLL can be extracted via PE-sieve 5](https://github.com/hasherezade/pe-sieve) or other tools for simplyfied further analysis.\n\nDepending on whether the previous process enumeration yielded results on the list, QakBot\nwill inject into `mobsync.exe (e.g., in case a` `avp.exe process is found running) instead of`\n```\nexplorer.exe . But for simplicity we will only follow the explorer.exe process injection\n\n```\npath we observed in our analysis environment.\n\n### C2 communication\n\nAfter avoiding detection, the injected QakBot code within `explorer.exe will start`\ncommunicating with the C2 servers.\n\n\n-----\n\nLike in previous versions of QakBot the C2 IP list is stored RC4 encrypted in resource\nsection `311 . The first 20 bytes of the section contains the RC4 key with which the rest of`\nthe section is decrypted. The first 20 bytes of the decrypted data will contain the SHA1 sum\ncalculated over the rest of the decrypted data. It is used as a verification for correct\ndecryption. Unlike in previous version, the C2 list is now stored in binary form and not as\nASCII text anymore.\n\nFor details on how to extract the C2 list and QakBot’s configuration see the Python3 script in\nthe appendix. The input to the script is the path to the DLL that QakBot injected into\n```\nexplorer.exe, which we previously extracted via PE-sieve .5\n\n```\nThe configuration is stored using the same RC4 encryption scheme in resource section\n```\n308 . In it we can see the bot and/or campaign ID abc103 that is associated with the\n\n```\nanalyzed sample. It is still stored in plain ACSII text. For each campaign the number is\n\n\n-----\n\nincreased by one. This allows the operators behind QakBot to keep track to which campaign\neach victim connecting to their C2 server belongs to. Another currently observed identifier is\n```\ntr02 . This identifier, however, stayed the same over multiple malspam campaigns.\n\n```\nVia the C2 connection the operators behind QakBot can remote control the malware and\ndeploy additional malicious modules.\n\nQakBot will not store its configuration and C2 list on disk anymore. It will use the registry for\nstorage.\n\n### Wiping\n\nThe previous QakBot version used to overwrite its initial executable with a copy of\n```\ncmd.exe . This version will overwrite the portion of the initially downloaded DLL after the PE\n\n```\nheader with zeros.\n\nHere is the entropy of the QakBot DLL as downloaded.\n\n\n-----\n\nThe zeroing of data after the header can be clearly seen when comparing the previous plot\nagainst a plot of the DLL file after wiping.\n\n### Persistence\n\nThe persistence mechanism of QakBot has also changed. While it still uses a run key\nregistry entry under `HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run, this key`\nis only set right before the system is shutdown, rebooted or put to sleep. The corresponding\nDLL is also only dropped to disk right before shutdown, rebooted or sleep.\n\nAfter the system boots up again, QakBot is started via the run key. The execution tree also\nstarts via `regsvr32.exe -s ... like the initial execution from Excel. QakBot follows the`\nsame steps as previously outlined resulting in process injection into `explorer.exe .`\n\nQakBot will then delete the run key registry entry and delete the DLL it dropped to disk prior\nto the reboot.\n\n\n-----\n\nThis way QakBot’s persistence can not be detected at runtime.\n\n### Egregor\n\n[While we have previously reported on QakBot deliverying the ProLock ransomware, 2](https://www.hornetsecurity.com/en/security-information/qakbot-malspam-leading-to-prolock/) latests\nreports indicated that QakBot is now used to deliver the Egregor ransomware. We previously\n[reported on the Egregor ransomware as part of an article on ransomware leaksites6](https://www.hornetsecurity.com/en/security-informationen-en/leakware-ransomware-hybrid-attacks/) in which\nwe explain the practice of ransomware operators stealing their victims data before encrypting\nit to extort them not only with decryption but also public release of the stolen data.\n\n\n## Conclusion and Countermeasures\n\nFrom our analysis we can conclude that QakBot is trying to avoid persistent file artifacts. In\nprevious version the configuration and QakBot executable were permanently stored on disk.\nThis made it easy for security tools to detect them. The new version tries to avoid\npermanently leaving its artifacts on disk. While QakBot is not going fully fileless, it new\ntactics will sure lower its detection.\n\nBut even though QakBot has changed, the delivery mechanism behind the QakBot “ abc[A```\nZ]+ ” campaign did not. Hence, an infection by this threat actor can be successfully\n\n```\nprevented by blocking the initial emails.\n\nHornetsecurity’s [Spam Filter and Malware Protection, with the highest detection rates on the](https://www.hornetsecurity.com/en/services/spam-filter/)\nmarket, already detects and blocks the outlined threat. Hornetsecurity’s Advanced Threat\nProtection extends this protection by also detecting yet unknown threats.\n\n## References\n\n Indicators of Compromise (IOCs)\n\n### Hashes\n\nThe hashes of the analyzed QakBot samples are:\n\n\n-----\n\n**MD5** **Filename** **Description**\n\n`6bc0584f6cbb74714add1718b0322655` `904400.jpg` QakBot DLL as downloaded\nby XLM macro\n\n`e23bc27212f61520cfb130185d74cfb1` `26e0000.dll` Extracted QakBot DLL\n\n## MITRE ATT&CK\n\nThe tactics and techniques used by QakBot as defined by the MITRE ATT&CK framework\nare as follows:\n\n**Tactic** **Technique**\n\nTA0001 – Initial Access T1566.001 – Phishing: Spearphishing Attachment\n\nTA0001 – Initial Access T1566.002 – Phishing: Spearphishing Link\n\nTA0002 – Execution T1027 – Obfuscated Files or Information\n\nTA0002 – Execution T1204.002 – User Execution: Malicious File\n\nTA0003 – Persistence T1547.001 – Boot or Logon Autostart Execution: Registry Run\nKeys / Startup Folder\n\n\nTA0004 – Privilege\nEscalation\n\nTA0005 – Defense\nEvasion\n\nTA0005 – Defense\nEvasion\n\nTA0005 – Defense\nEvasion\n\nTA0005 – Defense\nEvasion\n\nTA0005 – Defense\nEvasion\n\nTA0006 – Credential\nAccess\n\nTA0006 – Credential\nAccess\n\n\nT1053.005 – Scheduled Task/Job: Scheduled Task\n\nT1027.002 – Obfuscated Files or Information: Software Packing\n\nT1055 – Process Injection\n\nT1055.012 – Process Injection: Process Hollowing\n\nT1070 – Indicator Removal on Host\n\nT1497.001 – Virtualization/Sandbox Evasion: System Checks\n\nT1003 – OS Credential Dumping\n\nT1110.001 – Brute Force: Password Guessing\n\n\n-----\n\n**Tactic** **Technique**\n\n\nTA0006 – Credential\nAccess\n\nTA0011 – Command\nand Control\n\nTA0011 – Command\nand Control\n\nTA0011 – Command\nand Control\n\n## Appendix\n\n\nT1555.003 – Credentials from Password Stores: Credentials\nfrom Web Browsers\n\nT1071.001 – Application Layer Protocol: Web Protocols\n\nT1090 – Proxy\n\nT1090.002 – Proxy: External Proxy\n\n\n### Qakbot configuration extraction Python3 script\n```\nimport sys\nimport pefile\nfrom arc4 import ARC4\npe = pefile.PE(sys.argv[1])\nc2list = []\nfor entry in pe.DIRECTORY_ENTRY_RESOURCE.entries:\n  for e in entry.directory.entries:\n    n = e.name.string.decode()\n    data = pe.get_data(e.directory.entries[0].data.struct.OffsetToData,\ne.directory.entries[0].data.struct.Size)\n    data = ARC4(data[:20]).decrypt(data[20:])[20:]\n    if n == '311':\n      for i in range(1,len(data),7):\n        c2 = list(data[i:i+6])\n        c2list.append(\"%d.%d.%d.%d:%d\" % (c2[0],c2[1],c2[2],c2[3],(c2[4]\n<<8)+c2[5]))\n    elif n == '308':\n      config = data.decode().split()\nprint(\"# QakBot Config\\n\\n```\\n\" + \"\\n\".join(config) + \"\\n```\\n\")\nprint(\"# QakBot C2\\n\\n```\\n\" + \"\\n\".join(c2list) + \"\\n```\\n\")\n\n```\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2020/2020-12-15 - QakBot reducing its on disk artifacts.pdf"
    ],
    "report_names": [
        "2020-12-15 - QakBot reducing its on disk artifacts.pdf"
    ],
    "threat_actors": [],
    "ts_created_at": 1673535928,
    "ts_updated_at": 1743041127,
    "ts_creation_date": 1653713552,
    "ts_modification_date": 1653713552,
    "files": {
        "pdf": "https://archive.orkl.eu/7be67ab8cfee0db3022a05e06953a1169e202c8f.pdf",
        "text": "https://archive.orkl.eu/7be67ab8cfee0db3022a05e06953a1169e202c8f.txt",
        "img": "https://archive.orkl.eu/7be67ab8cfee0db3022a05e06953a1169e202c8f.jpg"
    }
}