# 0-raw-data/ORKL/library/jsons_entries/ef358fbb-1cfb-4e42-928a-e630248b63f7.json
{
    "id": "ef358fbb-1cfb-4e42-928a-e630248b63f7",
    "created_at": "2023-01-12T15:04:55.519596Z",
    "updated_at": "2025-03-27T02:05:21.049077Z",
    "deleted_at": null,
    "sha1_hash": "a20132dd2609c7b4a710f791e731da0673461e3d",
    "title": "2020-10-06 - The FONIX RaaS - New Low-Key Threat with Unnecessary Complexities",
    "authors": "",
    "file_creation_date": "2022-05-29T01:20:00Z",
    "file_modification_date": "2022-05-29T01:20:00Z",
    "file_size": 1407018,
    "plain_text": "# The FONIX RaaS | New Low-Key Threat with Unnecessary Complexities\n\n**[labs.sentinelone.com/the-fonix-raas-new-low-key-threat-with-unnecessary-complexities/](https://labs.sentinelone.com/the-fonix-raas-new-low-key-threat-with-unnecessary-complexities/)**\n\nJim Walter\n\n[FONIX Raas (Ransomware as a Service) is an offering that first came to attention in July of](https://www.sentinelone.com/wp-content/uploads/?s=RaaS)\nthis year. It did not make much of a splash at the time, and even currently, we are only\nseeing small numbers of infections due to this ransomware family. However, RaaS that at\nfirst fly under the radar can quickly become rampant if defenders and security solutions\nremain unaware of them. Notably, FONIX varies somewhat from many other current RaaS\nofferings in that it employs four methods of encryption for each file and has an overlycomplex post-infection engagement cycle. In this post, we dig a little deeper into these and\nother peculiarities of this new RaaS offering.\n\n## FONIX Background: From Crypters to Encrypters\n\nThe actors behind FONIX appeared to be primarily focused on binary crypters/packers prior\nto the release of the RaaS. Their ‘products’ were advertised on various cybercrime forums,\nas well paste-based advertisements on the Dark Web. Initial advertisement for the RaaS\nfollowed suit.\n\n\n-----\n\n## FONIX RaaS: A Complex Victim-Affiliate-Author Triangle\n\nEngagement for this RaaS is handled purely via email, and directly with the\nauthor/advertiser. There is no web-based portal to register or manage infections or\ncampaigns. The authors did appear to initially offer a FONIX-specific email service; however,\nat the time of writing, that service appears to be unavailable.\n\nUpon engaging with the FONIX advertiser, would-be buyers are required to supply the\nmalware author with their desired email address and password for the FONIX mail service.\nSince the FONIX mail service is currently inactive as noted, it appears that buyers are to\nsupply the sellers with alternative email addresses (e.g., protonmail). Once the seller has\nreceived the email data, the buyer is sent copies of the ransomware payloads.\n\nThe received payloads are customized to display the email address of the new buyer upon\ninfection, which in turn directs the victims to reach out via said email in order to receive\ndecryption instructions, or acquire proof of decryption. Again, all transactions are handled via\nemail, as opposed to a web-based portal.\n\nThere is no upfront cost for becoming a FONIX affiliate. Rather, when victims pay their\n[ransom (which they should ideally not do), the attacker (FONIX buyer) provides the FONIX](https://www.sentinelone.com/blog/ransomware-to-pay-or-not-to-pay-lets-discuss/)\nauthors with a 25% cut of the proceeds.\n\nThe actual process is a bit convoluted and far less user-friendly than most ransomware\nservices. Based on current intelligence, we know that FONIX affiliates do not get provided\nwith a decryptor utility or keys at first. Instead, victims first contact the affiliate (buyer) via\nemail as described above. The affiliate then requests a few files from the victim. These\ninclude two small files for decryption: one is to provide proof to the victim, the other is the file\n“cpriv.key” from the infected host. The affiliate is then required to send those files to the\nFONIX authors, who decrypt the files, after which they can be sent to the victims.\n\n\n-----\n\nPresumably, once the victim is satisfied that decryption is possible, the affiliate provides a\npayment address (BTC wallet). The victim then pays the affiliate, with the affiliate in turn\nsupplying the FONIX authors with their 25% cut.\n\nOnce the FONIX authors have received their portion of the proceeds, they provide the\naffiliate with the decryptor utility and key (keys are unique to each campaign). At that point, it\nis between the affiliate and the victim in terms of how they provide the decryption\ncapabilities.\n\nAll in all, this makes for a time-consuming process for any environment, especially large\nenterprises. Prevention of the infection, avoiding the whole rigmarole, is a far more attractive\noption!\n\n## FONIX Ransomware: File Encryption and Execution\n\nThe FONIX samples we have observed come in 64 and 32-bit varieties, and are available for\nWindows only. By default, FONIX will encrypt all file types, excluding critical Windows OS\nfiles.\n\nFile encryption is handled via a mixture of Salsa20, Chacha, RSA and AES.\n\nThe FONIX authors advertise that this is to ensure “strong” and “unbeatable” encryption.\nHowever, this does add considerable time to the encryption process. Our analysis shows that\nFONIX is between 2 and 5 times slower than other well-known ransomware families (e.g.,\nRyuk, [NetWalker).](https://www.sentinelone.com/wp-content/uploads/netwalker-ransomware-no-respite-no-english-required/)\n\nEncrypted files are all marked with the `.XINOF extension (FONIX backwards). Depending`\non the context of the executed payload, numerous other malicious changes are made to the\nsystem. In all cases, once encryption is complete, the Desktop background is changed to the\nFONIX logo, and the `.HTA -formatted ransomware note is displayed across the entire`\nscreen.\n\n\n-----\n\nAs noted, instructions to contact the attacker are provided in the ransom note ( How To\n```\nDecrypt Files.hta ). Several additional files are deposited on encrypted hosts. For\n\n```\nexample, the following can be found in `%programdata% post-encryption:`\n\nCpriv.key\nHello Michaele Gllips\n\n\n-----\n\nHelp.txt\nHow To Decrypt Files.hta\nSystemID\n```\nCpriv.key and SystemID are both required for decryption, as detailed in the decryption\n\n```\nchain process described above. `How To Decrypt Files.hta is the primary ransom note,`\nand this is the same HTA displayed prominently (covering up the FONIX logo wallpaper\nhowever). `Help.txt is an additional plain text file containing the same attacker email`\naddress. The file simply contains a quick message, the attacker email address and the\nSystemID.\n\nThe remaining file `Hello Michaele Gllips appears to be a message to @demonslay335`\nof the [MalwareHunterTeam. This was also documented, via Twitter, by](https://twitter.com/malwrhunterteam) [@bartblaze.](https://twitter.com/bartblaze/status/1291077314680369152)\n\nWhen executed with administrator privileges, the following additional system changes occur:\n\nTask Manager is disabled\nPersistence is achieved via scheduled task, Startup folder inclusion, and the registry\n(Run AND RunOnce)\nSystem file permissions are modified\nPersistent copies of the payload have their attributed set to hidden\nA hidden service is created for persistence (Windows 10)\nDrive / Volume labels are changed (to “XINOF”)\nVolume Shadow Copies are deleted (vssadmin, wmic)\nSystem recovery options are manipulated/disabled (bcdedit)\nSafeboot options are manipulated\n\n\n-----\n\nThe FONIX malware executes numerous system-altering commands. This is just a small\nsample of the many we observed:\n```\nconhost.exe 0xffffffff -ForceV1\ncmd.exe /c schtasks /CREATE /SC ONLOGON /TN fonix /TR C:ProgramDataXINOF.exe /RU\nSYSTEM /RL HIGHEST /F\nschtasks.exe /RU SYSTEM /RL HIGHEST /F\ncmd.exe /c copy XINOF.exe %appdata%MicrosoftWindowsStart MenuProgramsStartupXINOF.exe\ncmd.exe /c attrib +h +s C:Usersadmin1AppDataRoamingMicrosoftWindowsStart\nMenuProgramsStartupXINOF.exe\nattrib.exe\ncmd.exe /c reg add HKEY_LOCAL_MACHINESOFTWAREMicrosoftWindowsCurrentVersionRun /v\n\"PhoenixTechnology\" /t REG_SZ /d C:ProgramDataXINOF.exe /f\nreg.exe /f\ncmd.exe /c reg add HKEY_CURRENT_USERSOFTWAREMicrosoftWindowsCurrentVersionRun /v\n\"PhoenixTechnology\" /t REG_SZ /d C:ProgramDataXINOF.exe /f\ncmd.exe /c reg add HKEY_LOCAL_MACHINESOFTWAREMicrosoftWindowsCurrentVersionRunOnce /v\n\"PhoenixTechnology\" /t REG_SZ /d C:ProgramDataXINOF.exe /f\ncmd.exe /c reg add HKEY_CURRENT_USERSOFTWAREMicrosoftWindowsCurrentVersionRunOnce /v\n\"PhoenixTechnology\" /t REG_SZ /d C:ProgramDataXINOF.exe /f\ncmd.exe /c reg add\nHKEY_LOCAL_MACHINESoftwareMicrosoftWindowsCurrentVersionPoliciesSystem  /v\nDisableTaskMgr /t REG_DWORD /d 1 /f\nreg.exe add HKEY_LOCAL_MACHINESoftwareMicrosoftWindowsCurrentVersionPoliciesSystem \n/v DisableTaskMgr /t REG_DWORD /d 1 /f\ncmd.exe /c reg add\nHKEY_CURRENT_USERSoftwareMicrosoftWindowsCurrentVersionPoliciesSystem  /v\nDisableTaskMgr /t REG_DWORD /d 1 /f\n\n## Conclusion\n\n```\nFONIX is not at present a widespread threat; whether that is due to the complexity of its\nengagement model or other factors is difficult to say at this time. However, a FONIX infection\nis notably aggressive – encrypting everything other than system files – and can be difficult to\nrecover from once a device has been fully encrypted. Currently, FONIX does not appear to\nbe threatening victims with additional consequences (such as public data exposure or DDoS\nattacks) for non-compliance. Even without those extra headaches, however, good user\nhygiene and strong, modern, endpoint security controls are critical in preventing this and\nsimilar infections. The SentinelOne platform is fully able to prevent all behaviors and artifacts\nassociated with the FONIX ransomware family.\n\n\n-----\n\n## Indicators of Compromise\n\n**SHA1**\na94f92f1e6e4fed57ecb2f4ad55e22809197ba2e\n1f551246c5ed70e12371891f0fc6c2149d5fac6b\n63cae6a594535e8821c160da4b9a58fc71e46eb2\n\n**SHA256**\ne5324495a9328fe98187239565c05b077680b2ebc9183a6e3e2ccfbfa9f0295a\n5263c485f21886aad8737183a71ddc1dc77a92f64c58657c0628374e09bb6899\n658ec5aac2290606dba741bce30853515795028322162167395cebc5d0bfccf4\n\n**MITRE ATT&CK**\nData Encrypted for Impact [T1486](https://attack.mitre.org/techniques/T1486/)\n[Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder T1547]Obfuscated](https://attack.mitre.org/techniques/T1547/001/)\nFiles or Information [T1027](https://attack.mitre.org/techniques/T1027/)\n[Inhibit System Recovery T1490](https://attack.mitre.org/techniques/T1490/)\nScheduled Task/Job: Scheduled Task [T1053.005](https://attack.mitre.org/techniques/T1053/005/)\n[Boot or Logon Autostart Execution T1547](https://attack.mitre.org/techniques/T1547/)\n[Command and Scripting Interpreter T1059](https://attack.mitre.org/techniques/T1059/)\n[Command and Scripting Interpreter: Windows Command Shell T1059.003](https://attack.mitre.org/techniques/T1059/003/)\n[Obfuscated Files or Information: Software Packing T1027.002](https://attack.mitre.org/techniques/T1027/002/)\n[File and Directory Permissions Modification T1222](https://attack.mitre.org/techniques/T1222/)\nFile and Directory Permissions Modification: Windows File and Directory Permissions\n[Modification T1222.001](https://attack.mitre.org/techniques/T1222/001/)\nHide Artifacts [T1564](https://attack.mitre.org/techniques/T1564/)\n[Inhibit System Recovery T1490](https://attack.mitre.org/techniques/T1490/)\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2020/2020-10-06 - The FONIX RaaS - New Low-Key Threat with Unnecessary Complexities.pdf"
    ],
    "report_names": [
        "2020-10-06 - The FONIX RaaS - New Low-Key Threat with Unnecessary Complexities.pdf"
    ],
    "threat_actors": [],
    "ts_created_at": 1673535895,
    "ts_updated_at": 1743041121,
    "ts_creation_date": 1653787200,
    "ts_modification_date": 1653787200,
    "files": {
        "pdf": "https://archive.orkl.eu/a20132dd2609c7b4a710f791e731da0673461e3d.pdf",
        "text": "https://archive.orkl.eu/a20132dd2609c7b4a710f791e731da0673461e3d.txt",
        "img": "https://archive.orkl.eu/a20132dd2609c7b4a710f791e731da0673461e3d.jpg"
    }
}