# 0-raw-data/ORKL/library/jsons_entries/f4ec23e6-3591-4719-b253-6f99b8a6b2e1.json
{
    "id": "f4ec23e6-3591-4719-b253-6f99b8a6b2e1",
    "created_at": "2023-01-12T15:08:45.667485Z",
    "updated_at": "2025-03-27T02:06:14.529733Z",
    "deleted_at": null,
    "sha1_hash": "a216532ced7f312430a159260b23648897ebded5",
    "title": "2022-05-12 - Malware targeting latest F5 vulnerability",
    "authors": "",
    "file_creation_date": "2022-05-28T02:23:56Z",
    "file_modification_date": "2022-05-28T02:23:56Z",
    "file_size": 3109990,
    "plain_text": "# Malware targeting latest F5 vulnerability\n\n**[lacework.com/blog/malware-targeting-latest-f5-vulnerability/](https://www.lacework.com/blog/malware-targeting-latest-f5-vulnerability/)**\n\nChris Hall and Jared Stroud - Cloud Security Researchers, Lacework\n\nLacework Labs\n\nMay 12, 2022\n\n\nMay 12, 2022\n\n\n-----\n\nOn May 9th, 2022 a remote code execution vulnerability in F5’s [Big IP suite of appliances](https://www.f5.com/products/big-ip-services?&utm_medium=cpc&utm_source=google&utm_campaign=PSE_PG_22Q3_noa_big-ip-services_amer-sdla_iaa_ia&gclid=Cj0KCQjwmuiTBhDoARIsAPiv6L-K6cx6ruWV1ulm54z03UlOxccU019Zl9bnLONzC3TIFgB68vGXQuQaAncJEALw_wcB&gclsrc=aw.ds)\nunder [CVE-2022-1388. Per the CVE, appliance versions prior to 16.1.2.2, 15.1.x, 14.1.x,](https://cve.mitre.org/cgi-bin/cvename.cgi?name=CVE-2022-1388)\n14.1.4.6, 13.1.x, 13.1.5, and all versions within the 12.1.x and 11.6.x are affected. Notably,\nversions that have reached the end of technical support are not evaluated, and may also be\nvulnerable to this CVE. Since the announcement of this vulnerability numerous Github\nrepositories have been created showing proof-of-concept attacks which require nothing more\nthan a POST request with an HTTP body of commands to execute on a victim host.\nResearchers industry-wide have reported opportunistic adversaries adopting this\n[vulnerability to spread Cryptojacking (T1496), and DDoS bots (Mirai). Lacework Labs is also](https://attack.mitre.org/techniques/T1496/)\nbeginning to see payloads associated with this vulnerability within their honeypots. At this\ntime, Shodan reports 23,494 publicly facing BIG IP machines on the internet at the time of\nthis writing.\n\n\n-----\n\nPublic-Facing BIG-IP Appliances\n\nHunting for Malware\nLacework Labs has identified from their honeypot data CVE-2022-1388 being exploited.\nPost-execution activity ranges from simply executing “id”, to downloading and executing a\nsecond-stage payload. From this data, Lacework Labs developed YARA rules to hunt for the\nunique URI targeted in the F5 vulnerability and associated ELF files. This resulted in\nnumerous hits for Miria variants for various architectures demonstrating how quickly malware\nauthors can adopt PoCs to distribute their malware. When analyzing the executable\ncapabilities for exploiting other vulnerabilities aside from CVE-2022-1388 were discovered to\ninclude Log4J, ColdFusion, and various home router exploits. An example of this can be\nseen in the Ghidra pseudocode below highlighting a payload to take advantage of CVE2022-1388.\n\nFigure 1 – Example CVE-2022-1388 exploit template\n\nFurther examining the binary, a section of code with numerous XOR encoded strings were\nidentified. Given that the Mirai source code was leaked years ago, cross-referencing the\navailable source code gives further insight into changes an individual threat actor may have\nmade for their Mirai variant. The code snippet below shows the entries discovered within the\n[sample we analyzed on the left, and on the right shows the publicly available source code.](https://github.com/jgamblin/Mirai-Source-Code/blob/master/mirai/bot/scanner.c#L123)\n\nFigure 2 – Mirai source vs decompilation\n\nThe “add_entry” function contains a decoding routine that leverages a single key of “0x22”.\nThis allows us to easily decode said strings as shown in Figure 3. These strings are then\nused as credentials for brute force activity.\n\n\n-----\n\nFigure 3 – Decoding Strings\n\n[When executing, a crontab (T1053.003) entry is set to execute the binary every five minutes.](https://attack.mitre.org/techniques/T1053/003/)\n[Finally, during execution prctl is leveraged to rename the executing binary to a random set of](https://man7.org/linux/man-pages/man2/prctl.2.html)\ncharacters before continuing brute force efforts.\n\nFigure 4 – strace of\n\ncrontab entry\n\nFigure 5 – renaming binary pre-execution\nIdentifying Researchers\nAlong with exploitation in the wild using Mirai, Lacework Labs also observed activity\noriginating from Project Discovery. Project Discovery maintains a repo for a popular open\nsource scanner dubbed Nuclei which is used for both traditional scanning and recon as well\nas more invasive tests using OAST. To enable scanning the Nuclei community curates\nconfigs for various exploits. The configuration for CVE-2022-1388 was committed to Project\nDiscovery’s GitHub on May, 9th and Lacework has seen a moderate amount of related traffic\n[from varied sources. For more info on Project discovery check out our recent blog here.](https://www.lacework.com/blog/the-oast-with-the-most/)\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2022/2022-05-12 - Malware targeting latest F5 vulnerability.pdf"
    ],
    "report_names": [
        "2022-05-12 - Malware targeting latest F5 vulnerability.pdf"
    ],
    "threat_actors": [],
    "ts_created_at": 1673536125,
    "ts_updated_at": 1743041174,
    "ts_creation_date": 1653704636,
    "ts_modification_date": 1653704636,
    "files": {
        "pdf": "https://archive.orkl.eu/a216532ced7f312430a159260b23648897ebded5.pdf",
        "text": "https://archive.orkl.eu/a216532ced7f312430a159260b23648897ebded5.txt",
        "img": "https://archive.orkl.eu/a216532ced7f312430a159260b23648897ebded5.jpg"
    }
}