# 0-raw-data/ORKL/library/jsons_entries/ffc9fd48-3ea3-46f5-ae38-027460d7d11b.json
{
    "id": "ffc9fd48-3ea3-46f5-ae38-027460d7d11b",
    "created_at": "2023-01-12T15:05:53.870351Z",
    "updated_at": "2025-03-27T02:05:23.813391Z",
    "deleted_at": null,
    "sha1_hash": "938766693a3ef528f13cb909b392e96ca46aca6c",
    "title": "2021-10-28 - Spook Ransomware - Prometheus Derivative Names Those That Pay, Shames Those That Don’t",
    "authors": "",
    "file_creation_date": "2022-05-29T01:20:17Z",
    "file_modification_date": "2022-05-29T01:20:17Z",
    "file_size": 2016120,
    "plain_text": "# Spook Ransomware | Prometheus Derivative Names Those That Pay, Shames Those That Don’t\n\n**[sentinelone.com/labs/spook-ransomware-prometheus-derivative-names-those-that-pay-shames-those-that-dont/](https://www.sentinelone.com/labs/spook-ransomware-prometheus-derivative-names-those-that-pay-shames-those-that-dont/)**\n\nJim Walter\n\n**_By Jim Walter and Niranjan Jayanand_**\n\n## Executive Summary\n\nSpook Ransomware is an emerging player first seen in late September 2021\nThe operators publish details of all victims regardless of whether they pay or not\nTargets range across several industries with an emphasis on manufacturing\nAnalysis shows a significant degree of code sharing between Spook and the\nPrometheus and Thanos ransomware families\n\n## Overview\n\nSpook ransomware emerged onto the scene in late September 2021 and follows the multipronged extortion model that is all too common these days. Victims are hit with the threat of\ndata destruction as well as public data leakage and the associated fallout. In this report, we\nexplore how the malware shares certain similarities with earlier ransomware families, and\ndescribe its main encryption and execution behaviour.\n\n## Spook and Prometheus\n\n\n-----\n\n[There is some indication that Spook is either linked to, or derived from, Prometheus](https://www.youtube.com/watch?v=CgyxMQFFU_k)\n[ransomware. Prometheus is itself an evolution of Thanos ransomware. However, it is](https://www.sentinelone.com/labs/thanos-ransomware-riplace-bootlocker-and-more-added-to-feature-set/)\nimportant to note that since Thanos ransomware had a builder which was leaked, any real\nattempts at attribution based solely on the malware’s code is somewhat futile. Even so, there\nare a few notable similarities between Spook, Prometheus, and ultimately Thanos.\n\nThe .NET binary in the following sample, first seen in VirusTotal on 02 October, provides a\nglimpse into some of these similarities, with artifacts from the Thanos builder also apparent.\n\n```\na63a5de26582af1438c9886cfb15c4baa08cce2e\n\n```\n\nShared code block with Thanos\n\nOur analysis suggests that there is an overlap of between 29-50% of shared code between\nSpook and Prometheus. Some of this overlap is related to construction of the ransom notes\nand key identifiers.\n\n\n-----\n\nRansom note similarity example (Prometheus vs Spook)\nIn addition to shared code artifacts, there are similarities with regards to the layout and\nstructure of the Spook and Prometheus payment portals.\n\nBelow are the similarities between the leak data URLs hosted by both the groups\n\nSpook ransomware:\n```\n   hxxp[:]//spookuhv****.onion/blog/wp-content/uploads/2021/05/1-15.png\n\n```\nPrometheus ransomware:\n\n```\nhxxp[:]//promethw****.onion/blog/wp-content/uploads/2021/05/1-15.png\n\n```\n\n## Offline Encryption and Process Manipulation\n\nSpook, mirroring the manifestos of others, boasts “very strong (AES) encryption” along with\nthe threat of leaking victim data to the public. The malware has the ability to encrypt target\nmachines without requiring internet connectivity. Encryption of a full disk can occur within just\na few minutes, at which point the ransom note is displayed on the desktop\n( RESTORE_FILES_INFO.HTA ) along with numerous other system notifications.\n\nThe malware also makes a number of changes to ensure that the ransom notifications are\ndisplayed prominently after reboot (via Start Menu lnk, Reg).\n```\nWinLogon is modified (via registry) to display the Ransom Note text upon login:\n     HKLM\\SOFTWARE\\Wow6432Node\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\n     Str Value: LegalNoticeCaption/Text\n\n```\n\n-----\n\nRegistry Modifications for Persistence\nRansom notes are also displayed upon login via a Shortcut placed in the Startup directory\n\n\n-----\n\nStartup Folder Shortcut\nIn addition, Spook will attempt to terminate processes and stop services of anything that may\ninhibit the encryption process.\n\nHere again there is overlap between Spook, Prometheus, and Thanos with regards to\nprocess discovery and manipulation, especially with regards to checking for and killing the\n[Raccine anti-ransomware process that some organizations deploy in an effort to protect](https://github.com/Neo23x0/Raccine#readme)\nshadow copies.\n```\nTASKILL.EXE is used to force the termination of the following processes if found:\n\n```\n\n-----\n\n```\n     agntsvc.exe\n     CNTAoSMgr.exe\n     dbeng50.exe\n     dbsnmp.exe\n     encsvc.exe\n     excel.exe\n     firefoxconfig.exe\n     hunderbird.exe\n     infopath.exe\n     isqlplussvc.exe\n     mbamtray.exe\n     msaccess.exe\n     msftesql.exe\n     mydesktopqos.exe\n     mydesktopservice.exe\n     mysqld-nt.exe\n     Mysqld-opt.exe\n     Mspub.exe\n     mysqld.exe\n     Ntrtscan.exe\n     ocautoupds.exe\n     ocomm.exe\n     ocssd.exe\n     onenote.exe\n     oracle.exe\n     outlook.exe\n     PccNTMon.exe\n     Powerpnt.exe\n     RaccineSettings.exe\n     sqbcoreservice.exe\n     sqlagent.exe\n     sqlbrowser.exe\n     sqlservr.exe\n     Sqlwriter.exe\n     synctime.exe\n     steam.exe\n     tbirdconfig.exe\n     thebat.exe\n     thebat64.exe\n     tmlisten.exe\n     visio.exe\n     winword.exe\n     wordpad.exe\n     xfssvccon.exe\n     zoolz.exe\n     taskkill.exe /IM ocomm.exe /F\n\n```\nThe Raccine product is specifically targeted with regards to disabling the products’ UI\ncomponents and update features. These are carried out via basic OS commands such as\n```\nreg.exe and schtasks.exe .\n\n```\n\n-----\n\n```\n     taskkill.exe /F /IM RaccineSettings.exe\n     reg.exe (CLI interpreter) delete\n\"HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\" /V \"Raccine Tray\" /F\n     reg.exe (CLI interpreter) delete HKCU\\Software\\Raccine /F\n     schtasks.exe (CLI interpreter) /DELETE /TN \"Raccine Rules Updater\" /F\n\n```\nIn addition, `sc.exe is used to disable specific services and components:`\n```\n     sc.exe config Dnscache start= auto\n     sc.exe config SQLTELEMETRY start= disabled\n     sc.exe config FDResPub start= auto\n     sc.exe config SSDPSRV start= auto\n     sc.exe config SQLTELEMETRY$ECWDB2 start= disabled\n     sc.exe config SstpSvc start= disabled\n     sc.exe config upnphost start= auto\n     sc.exe config SQLWriter start= disabled\n\n```\nWith various processes out of the way and the system in an optimal state for encryption, the\nmalware proceeds to enumerate local files and folders, along with accessible network\nresources.\n\nGiven the Thanos pedigree, specifics around encryption can vary. The samples analyzed\nemploy a random string at runtime as the passphrase for file encryption (AES). The string is\nsubsequently encrypted with the attacker’s public key and added into the generated ransom\nnote(s). Recovery of encrypted data is, therefore, not possible without the corresponding\nprivate key.\n\n## Ransom Payment and Victimology\n\nUpon infection, victims are instructed to proceed to Spook’s TOR-based payment portal.\n\n\n-----\n\nSpook Ransom Demand\nAt the payment portal, the victim is able to interact with the attackers via chat to negotiate\npayment.\n\n\n-----\n\nSpook Payment Portal\nSpook has been leveraging attacks against high-value targets across the globe, with little to\nno discretion with regards to industry. Looking at the current cross-section of victims posted\non the group’s web site, however, the majority are in the manufacturing sector.\n\nThe public blog went live in early October 2021. At the time of writing, there are 17 victims\nposted on the Spook site.\n\n\n-----\n\nSome of the victims named on the Spook blog site\nSpook actually lists all attacked companies, regardless of whether or not they pay the\nransom demand. Those victims that pay have their entry updated to indicate that the\ncompany’s data is ‘not for sale’. Those that have not paid are listed as having data that is\n“For Sale”, while some victim entries, presumably the most recent or those that are in the\nprocess of negotiating, are listed as “Company Decides”.\n\n\n-----\n\n## Conclusion\n\nAs these attacks continue to escalate and become more egregious, the need for true attack\nprevention is all the more critical. Spook’s tactic of public outing victims even if they pay\nthreatens reputational harm to any compromised company, even if they follow the attackers’\npayment demands.\n\nThis only continues to illustrate the importance of preventing attacks in the first place.\nRansomware operators have moved beyond worrying about companies detecting after-thefact and attempting to recover encrypted data.\n\n## Indicators of Compromise\n\n**SHA256**\n\n8dad29bd09870ab9cacfdea9e7ab100d217ff128aea64fa4cac752362459991c\n\ne347fd231a543a5dfd53b01ff0bc67b2bf37593e7ddc036f15bac8ad92f0d707\n\n\n-----\n\nd991aa2b1fad608b567be28e2d13d3d4f48eea3dea8f5d51a8e42aa9a2637426\n\n**SHA1**\na63a5de26582af1438c9886cfb15c4baa08cce2e\nbfd0ab7eec4b282cc5689a48e8f438d042c9d98f\ne2b098d36e51d2b7405fadbd578cf9774433f85a\n\n**MITRE ATT&CK**\n[TA0005 – Defense Evasion](https://attack.mitre.org/tactics/TA0005/)\n[T1486 – Data Encrypted for Impact](https://attack.mitre.org/techniques/T1486/)\n[T1027.002 – Obfuscated Files or Information: Software Packing](https://attack.mitre.org/techniques/T1027/002/)\n[T1007 – System Service Discovery](https://attack.mitre.org/techniques/T1007/)\n[T1059 – Command and Scripting Interpreter](https://attack.mitre.org/techniques/T1059/)\n[T1112 – Modify Registry](https://attack.mitre.org/techniques/T1112/)\n[TA0010 – Exfiltration](https://attack.mitre.org/tactics/TA0010/)\n[T1018 – Remote System Discovery](https://attack.mitre.org/techniques/T1018/)\n[T1082 – System Information Discovery](https://attack.mitre.org/techniques/T1082/)\n[T1547.004 – Boot or Logon Autostart Execution: Winlogon Helper DLL](https://attack.mitre.org/techniques/T1547/004/)\n[T1547.001 – Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder](https://attack.mitre.org/techniques/T1547/001/)\n\n\n-----\n\n**Spook Ransom Note Sample**\n\n\n-----",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2021/2021-10-28 - Spook Ransomware - Prometheus Derivative Names Those That Pay, Shames Those That Don’t.pdf"
    ],
    "report_names": [
        "2021-10-28 - Spook Ransomware - Prometheus Derivative Names Those That Pay, Shames Those That Don’t.pdf"
    ],
    "threat_actors": [
        {
            "id": "c28760b2-5ec6-42ad-852f-be00372a7ce4",
            "created_at": "2022-10-27T08:27:13.172734Z",
            "updated_at": "2025-03-27T02:00:55.511386Z",
            "deleted_at": null,
            "main_name": "Ember Bear",
            "aliases": [
                "Ember Bear",
                "UNC2589",
                "Bleeding Bear",
                "DEV-0586",
                "Cadet Blizzard",
                "Frozenvista",
                "UAC-0056"
            ],
            "source_name": "MITRE:Ember Bear",
            "tools": [
                "P.A.S. Webshell",
                "CrackMapExec",
                "ngrok",
                "WhisperGate",
                "Saint Bot",
                "PsExec",
                "Rclone",
                "Impacket",
                "OutSteel"
            ],
            "source_id": "MITRE",
            "reports": null
        },
        {
            "id": "083d63b2-3eee-42a8-b1bd-54e657a229e8",
            "created_at": "2022-10-25T16:07:24.143338Z",
            "updated_at": "2025-03-27T02:02:10.12138Z",
            "deleted_at": null,
            "main_name": "SaintBear",
            "aliases": [
                "Ember Bear",
                "FROZENVISTA",
                "Lorec53",
                "Nascent Ursa",
                "Nodaria",
                "SaintBear",
                "Storm-0587",
                "TA471",
                "UAC-0056",
                "UNC2589"
            ],
            "source_name": "ETDA:SaintBear",
            "tools": [
                "Agentemis",
                "Cobalt Strike",
                "CobaltStrike",
                "Elephant Client",
                "Elephant Implant",
                "GraphSteel",
                "Graphiron",
                "GrimPlant",
                "OutSteel",
                "Saint Bot",
                "SaintBot",
                "cobeacon"
            ],
            "source_id": "ETDA",
            "reports": null
        },
        {
            "id": "eecf54a2-2deb-41e5-9857-fed94a53f858",
            "created_at": "2023-01-06T13:46:39.349959Z",
            "updated_at": "2025-03-27T02:00:03.057647Z",
            "deleted_at": null,
            "main_name": "SaintBear",
            "aliases": [
                "UAC-0056",
                "Nascent Ursa",
                "Nodaria",
                "FROZENVISTA",
                "Storm-0587",
                "Lorec53",
                "Lorec Bear",
                "UNC2589",
                "TA471",
                "DEV-0587",
                "Saint Bear",
                "EMBER BEAR",
                "Bleeding Bear"
            ],
            "source_name": "MISPGALAXY:SaintBear",
            "tools": [],
            "source_id": "MISPGALAXY",
            "reports": null
        }
    ],
    "ts_created_at": 1673535953,
    "ts_updated_at": 1743041123,
    "ts_creation_date": 1653787217,
    "ts_modification_date": 1653787217,
    "files": {
        "pdf": "https://archive.orkl.eu/938766693a3ef528f13cb909b392e96ca46aca6c.pdf",
        "text": "https://archive.orkl.eu/938766693a3ef528f13cb909b392e96ca46aca6c.txt",
        "img": "https://archive.orkl.eu/938766693a3ef528f13cb909b392e96ca46aca6c.jpg"
    }
}